# beta/macOS/System/Library/LaunchAgents/com.apple.appleaccountd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.appleaccountd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.jett.switch.environmentChange.idms.complete</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.jett.switch.environmentChange.idms.complete</string>
			</dict>
			<key>com.apple.keystore.lockstatus</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.keystore.lockstatus</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.appleaccountd.appleid.availability</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.appleaccountd.postinstall</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<real>3600</real>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.appleaccountd.trustedContacts.healthCheck</key>
			<dict>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aa.accountService.xpc</key>
		<true/>
		<key>com.apple.aa.custodian.xpc</key>
		<true/>
		<key>com.apple.aa.daemon.xpc</key>
		<true/>
		<key>com.apple.aa.identity.xpc</key>
		<true/>
		<key>com.apple.aa.inheritance.xpc</key>
		<true/>
		<key>com.apple.aa.usernotifications.xpc</key>
		<true/>
		<key>com.apple.aps.aa</key>
		<true/>
		<key>com.apple.aps.appleaccountd</key>
		<true/>
		<key>com.apple.private.alloy.accounts.representative-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.appleaccount.notifications</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/appleaccountd</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>3</integer>
</dict>
</plist>
