# beta/macOS/System/Library/LaunchAgents/com.apple.searchtoold.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>HardResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>1024</integer>
	</dict>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamMemoryLimit</key>
		<integer>4000</integer>
		<key>JetsamPriority</key>
		<integer>-49</integer>
	</dict>
	<key>Label</key>
	<string>com.apple.searchtoold</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.omniSearch.answerSynthesis</key>
		<true/>
		<key>com.apple.omniSearch.echo</key>
		<true/>
		<key>com.apple.omniSearch.personalSearch</key>
		<true/>
		<key>com.apple.omniSearch.search</key>
		<true/>
	</dict>
	<key>NoEnvironmentVariables</key>
	<dict>
		<key>MallocStackLogging</key>
		<string>1</string>
		<key>no_DYLD_INSERT_LIBRARIES</key>
		<string>/usr/lib/libgmalloc.dylib</string>
		<key>no_MallocStackLoggingNoCompact</key>
		<string>1</string>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/OmniSearch.framework/Versions/A/searchtoold</string>
	</array>
	<key>SoftResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>1024</integer>
	</dict>
	<key>_ExponentialThrottling</key>
	<true/>
</dict>
</plist>
