# beta/macOS/System/Library/LaunchAgents/com.apple.mediaanalysisd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.mediaanalysisd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.mediaanalysisd.background.scheduler</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.mediaanalysisd</string>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>3600</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>2900</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresUserInactivity</key>
				<true/>
				<key>ResourceIntensive</key>
				<false/>
			</dict>
			<key>com.apple.mediaanalysisd.preheat</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.mediaanalysisd</string>
				<key>NonRepeatingTask</key>
				<dict>
					<key>ScheduleAfter</key>
					<integer>0</integer>
					<key>TrySchedulingBefore</key>
					<integer>60</integer>
				</dict>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresUserInactivity</key>
				<false/>
				<key>ResourceIntensive</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mediaanalysisd.analysis</key>
		<true/>
		<key>com.apple.mediaanalysisd.computeservice</key>
		<true/>
		<key>com.apple.mediaanalysisd.embeddingstore</key>
		<true/>
		<key>com.apple.mediaanalysisd.homekit</key>
		<true/>
		<key>com.apple.mediaanalysisd.homekitsession</key>
		<true/>
		<key>com.apple.mediaanalysisd.photos</key>
		<true/>
		<key>com.apple.mediaanalysisd.realtime</key>
		<true/>
		<key>com.apple.mediaanalysisd.service.public</key>
		<true/>
		<key>com.apple.mediaanalysisd.videosession.public</key>
		<true/>
		<key>com.apple.mediaanalysisd.xpcstore.spl.vuindex</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/MediaAnalysis.framework/Versions/A/mediaanalysisd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
