# beta/macOS/System/Library/LaunchAgents/com.apple.coreservices.useractivityd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>DrainMessagesAfterFailedInit</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.coreservices.useractivityd</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
		<string>LoginWindow</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.coreservices.lsbestappsuggestionmanager.xpc</key>
		<true/>
		<key>com.apple.coreservices.lsuseractivity.simulatorsupport</key>
		<true/>
		<key>com.apple.coreservices.lsuseractivityadministration.xpc</key>
		<true/>
		<key>com.apple.coreservices.lsuseractivitymanager.xpc</key>
		<true/>
		<key>com.apple.coreservices.uasharedpasteboardaux.xpc</key>
		<true/>
		<key>com.apple.coreservices.uasharedpasteboardcontroll.xpc</key>
		<true/>
		<key>com.apple.coreservices.uasharedpasteboardmanager.xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/UserActivity.framework/Agents/useractivityd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
