# beta/macOS/System/Library/LaunchAgents/com.apple.systemsettingsagent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.systemsettingsagent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>blockStorage-device-blueray.connected</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IOBDBlockStorageDevice</string>
			</dict>
			<key>blockStorage-device-blueray.disconnected</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IONotificationType</key>
				<string>IOServiceTerminate</string>
				<key>IOProviderClass</key>
				<string>IOBDBlockStorageDevice</string>
			</dict>
			<key>blockStorage-device-cd.connected</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IOCDBlockStorageDevice</string>
			</dict>
			<key>blockStorage-device-cd.disconnected</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IONotificationType</key>
				<string>IOServiceTerminate</string>
				<key>IOProviderClass</key>
				<string>IOCDBlockStorageDevice</string>
			</dict>
			<key>blockStorage-device-dvd.connected</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IODVDBlockStorageDevice</string>
			</dict>
			<key>blockStorage-device-dvd.disconnected</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IONotificationType</key>
				<string>IOServiceTerminate</string>
				<key>IOProviderClass</key>
				<string>IODVDBlockStorageDevice</string>
			</dict>
			<key>keyboard.connected</key>
			<dict>
				<key>DeviceUsagePairs</key>
				<array>
					<dict>
						<key>DeviceUsage</key>
						<integer>6</integer>
						<key>DeviceUsagePage</key>
						<integer>1</integer>
					</dict>
				</array>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IOHIDDevice</string>
			</dict>
			<key>keyboard.disconnected</key>
			<dict>
				<key>DeviceUsagePairs</key>
				<array>
					<dict>
						<key>DeviceUsage</key>
						<integer>6</integer>
						<key>DeviceUsagePage</key>
						<integer>1</integer>
					</dict>
				</array>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IONotificationType</key>
				<string>IOServiceTerminate</string>
				<key>IOProviderClass</key>
				<string>IOHIDDevice</string>
			</dict>
			<key>pointing-device.connected</key>
			<dict>
				<key>DeviceUsagePairs</key>
				<array>
					<dict>
						<key>DeviceUsage</key>
						<integer>1</integer>
						<key>DeviceUsagePage</key>
						<integer>1</integer>
					</dict>
				</array>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IONotificationType</key>
				<string>IOServicePublish</string>
				<key>IOProviderClass</key>
				<string>IOHIDEventService</string>
			</dict>
			<key>pointing-device.disconnected</key>
			<dict>
				<key>DeviceUsagePairs</key>
				<array>
					<dict>
						<key>DeviceUsage</key>
						<integer>1</integer>
						<key>DeviceUsagePage</key>
						<integer>1</integer>
					</dict>
				</array>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IONotificationType</key>
				<string>IOServiceTerminate</string>
				<key>IOProviderClass</key>
				<string>IOHIDDevice</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.BiometricKit.connectedAccessoriesChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.BiometricKit.connectedAccessoriesChanged</string>
			</dict>
			<key>com.apple.systemsettings.extensions.didchange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.systemsettings.extensions.didchange</string>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Applications/System Settings.app/Contents/Resources/systemsettingsagent</string>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
