# beta/macOS/System/Library/LaunchAgents/com.apple.triald.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.triald</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.accounts.changes</key>
		<dict>
			<key>AppleAccountAddedOrDeleted</key>
			<dict>
				<key>com.apple.account.AppleAccount</key>
				<dict>
					<key>ChangeTypes</key>
					<array>
						<integer>1</integer>
						<integer>3</integer>
					</array>
				</dict>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AFDictationEnablementDidChangeDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>AFDictationEnablementDidChangeDarwinNotification</string>
			</dict>
			<key>AFLanguageCodeDidChangeDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>AFLanguageCodeDidChangeDarwinNotification</string>
			</dict>
			<key>AppleKeyboardsSettingsChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleKeyboardsSettingsChangedNotification</string>
			</dict>
			<key>SUPreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>SUPreferencesChangedNotification</string>
			</dict>
			<key>com.apple.GeoServices.countryCodeChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.GeoServices.countryCodeChanged</string>
			</dict>
			<key>com.apple.GeoServices.experimentsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.GeoServices.experimentsChanged</string>
			</dict>
			<key>com.apple.ManagedConfiguration.profileListChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ManagedConfiguration.profileListChanged</string>
			</dict>
			<key>com.apple.SoftwareUpdate.SUPreferencesChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SoftwareUpdate.SUPreferencesChanged</string>
			</dict>
			<key>com.apple.crashreporter.auto_submit_preference_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.crashreporter.auto_submit_preference_changed</string>
			</dict>
			<key>com.apple.gms.availability.notification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.gms.availability.notification</string>
			</dict>
			<key>com.apple.language.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.language.changed</string>
			</dict>
			<key>com.apple.managedconfiguration.effectivesettingschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.effectivesettingschanged</string>
			</dict>
			<key>com.apple.triald.new-experiment</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.triald.new-experiment</string>
			</dict>
			<key>com.apple.triald.wake</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.triald.wake</string>
			</dict>
			<key>kAFPreferencesDidChangeDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>kAFPreferencesDidChangeDarwinNotification</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.triald.cellular</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>259200</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<integer>3600</integer>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.triald.client-trigger.cellular</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<false/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.triald.client-trigger.wifi</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<false/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.triald.fetch-experiments</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<integer>3600</integer>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.triald.maintenance</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<real>3600</real>
				<key>Repeating</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.triald.post-upgrade</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.triald.post-upgrade-require-inexpensive-networking</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>21600</integer>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.triald.retries-networking</key>
			<dict/>
			<key>com.apple.triald.setup-assistant-fetch</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>DataBudgeting</key>
				<true/>
				<key>GracePeriod</key>
				<integer>21600</integer>
				<key>NetworkDownloadSize</key>
				<integer>2147483648</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.triald.task-queue</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.triald</key>
		<true/>
		<key>com.apple.trial.status</key>
		<true/>
		<key>com.apple.triald.cache-delete</key>
		<true/>
		<key>com.apple.triald.internal</key>
		<true/>
		<key>com.apple.triald.namespace-management</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/triald</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/triald</string>
	</array>
</dict>
</plist>
