# beta/macOS/System/Library/LaunchAgents/com.apple.cfprefsd.xpc.agent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>DrainMessagesOnFailedInit</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>HardResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>512</integer>
	</dict>
	<key>Label</key>
	<string>com.apple.cfprefsd.xpc.agent</string>
	<key>LimitLoadToSessionType</key>
	<string>Background</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.cfprefsd.agent</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/cfprefsd</string>
		<string>agent</string>
	</array>
	<key>SoftResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>512</integer>
	</dict>
</dict>
</plist>
