# beta/macOS/System/Library/LaunchAgents/com.apple.mobiletimerd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>OS_ACTIVITY_ENABLE_DYNAMIC</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.mobiletimerd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.alarm</key>
		<true/>
		<key>com.apple.coreduetcontext.client_event_stream</key>
		<dict>
			<key>notification</key>
			<dict>
				<key>client-identifier</key>
				<string>com.apple.mobiletimerd.contextstore-registration</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>Application Uninstalled</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>SignificantTimeChange</key>
			<dict>
				<key>Notification</key>
				<string>SignificantTimeChangeNotification</string>
			</dict>
			<key>com.apple.Carousel.wristStateChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.Carousel.wristStateChanged</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.mobiletimerd.bedtimetest</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobiletimerd.bedtimetest</string>
			</dict>
			<key>com.apple.mobiletimerd.chargetest</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobiletimerd.chargetest</string>
			</dict>
			<key>com.apple.mobiletimerd.diagnostics</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobiletimerd.diagnostics</string>
			</dict>
			<key>com.apple.mobiletimerd.goodmorningtest</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobiletimerd.goodmorningtest</string>
			</dict>
			<key>com.apple.mobiletimerd.reset</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobiletimerd.reset</string>
			</dict>
			<key>com.apple.mobiletimerd.resttest</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobiletimerd.resttest</string>
			</dict>
			<key>com.apple.mobiletimerd.waketest</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobiletimerd.waketest</string>
			</dict>
			<key>com.apple.mobiletimerd.wakeuptest</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobiletimerd.wakeuptest</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidunpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidunpair</string>
			</dict>
			<key>com.apple.nanoregistry.paireddevicedidchangecapabilities</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.paireddevicedidchangecapabilities</string>
			</dict>
			<key>com.apple.nanoregistry.paireddevicedidchangeversion</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.paireddevicedidchangeversion</string>
			</dict>
			<key>com.apple.nanoregistry.watchdidbecomeactive</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.watchdidbecomeactive</string>
			</dict>
			<key>com.apple.springboard.finishedstartup</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.finishedstartup</string>
			</dict>
			<key>com.apple.system.powersources.source</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powersources.source</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.MobileTimer.alarmserver</key>
		<true/>
		<key>com.apple.MobileTimer.sessionserver</key>
		<true/>
		<key>com.apple.MobileTimer.stopwatchserver</key>
		<true/>
		<key>com.apple.MobileTimer.timerserver</key>
		<true/>
		<key>com.apple.alarmkitservices</key>
		<true/>
		<key>com.apple.mobiletimer.aps</key>
		<true/>
		<key>com.apple.private.alloy.mobiletimersync-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.clock</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.mobiletimer</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/MobileTimer.framework/Executables/mobiletimerd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
