# beta/macOS/System/Library/LaunchAgents/com.apple.cloudd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.cloudd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.cdp.manatee-availability</key>
			<dict>
				<key>Name</key>
				<string>com.apple.cdp.manatee-availability</string>
			</dict>
			<key>com.apple.cdp.manatee-availability-changed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.cdp.manatee-availability-changed</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.cloudasset.cloudd</key>
		<true/>
		<key>com.apple.cloudd</key>
		<true/>
		<key>com.apple.cloudd.cache-delete</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CloudKitDaemon.framework/Support/cloudd</string>
	</array>
</dict>
</plist>
