# beta/macOS/System/Library/LaunchAgents/com.apple.assistantd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.assistantd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.siri.bg_system_task.speech-audio-sample-delete</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.siri.bg_system_task.speech-audio-sampled</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>ABDatabaseChangedExternallyNotificationPriv</key>
			<dict>
				<key>Name</key>
				<string>ABDatabaseChangedExternallyNotificationPriv</string>
			</dict>
			<key>ABDistributedDatabaseChangedNotificationPriv</key>
			<dict>
				<key>Name</key>
				<string>ABDistributedDatabaseChangedNotificationPriv</string>
			</dict>
			<key>CKPersistentStorageKeysAndValuesUpdated.SiriUserCorrections</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysAndValuesUpdated.SiriUserCorrections</string>
			</dict>
			<key>CKPersistentStorageKeysAndValuesUpdated.synched</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysAndValuesUpdated.synched</string>
			</dict>
			<key>CKPersistentStorageKeysRemoved.SiriUserCorrections</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysRemoved.SiriUserCorrections</string>
			</dict>
			<key>CKPersistentStorageKeysRemoved.synched</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysRemoved.synched</string>
			</dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AppleDatePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleDatePreferencesChangedNotification</string>
			</dict>
			<key>MPStoreClientTokenDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>MPStoreClientTokenDidChangeNotification</string>
			</dict>
			<key>com.apple.MobileAsset.AutoAssetNotification^com.apple.MobileAsset.MAAutoAsset^STARTUP_ACTIVATED</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.AutoAssetNotification^com.apple.MobileAsset.MAAutoAsset^STARTUP_ACTIVATED</string>
			</dict>
			<key>com.apple.MobileAsset.EmbeddedSpeechMac.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.EmbeddedSpeechMac.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceServices.VoiceResources.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceServices.VoiceResources.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsASMac.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsASMac.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsASMac.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsASMac.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsMac.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsMac.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsMac.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsMac.new-asset-installed</string>
			</dict>
			<key>com.apple.SafariShared.Assistant.reload_plugin</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SafariShared.Assistant.reload_plugin</string>
			</dict>
			<key>com.apple.VideoSubscriberAccount.DidRegisterSubscription</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.VideoSubscriberAccount.DidRegisterSubscription</string>
			</dict>
			<key>com.apple.assistant.app_vocabulary</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.app_vocabulary</string>
			</dict>
			<key>com.apple.assistant.siri_settings_did_change</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.siri_settings_did_change</string>
			</dict>
			<key>com.apple.assistant.sync_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.sync_data_changed</string>
			</dict>
			<key>com.apple.assistant.sync_homekit_now</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.sync_homekit_now</string>
			</dict>
			<key>com.apple.assistant.sync_homekit_urgent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.sync_homekit_urgent</string>
			</dict>
			<key>com.apple.assistant.sync_needed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.sync_needed</string>
			</dict>
			<key>com.apple.chatkit.groups.siri_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.chatkit.groups.siri_data_changed</string>
			</dict>
			<key>com.apple.geoservices.siri_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.geoservices.siri_data_changed</string>
			</dict>
			<key>com.apple.icloud.fmip.siri_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.fmip.siri_data_changed</string>
			</dict>
			<key>com.apple.managedconfiguration.effectivesettingschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.effectivesettingschanged</string>
			</dict>
			<key>com.apple.siri.client.state.DynamiteClientState.siri_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.client.state.DynamiteClientState.siri_data_changed</string>
			</dict>
			<key>com.apple.siri.preheat.quiet</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.preheat.quiet</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.siri.understanding</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.siri.understanding</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.siri.understanding.nl.overrides</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.siri.understanding.nl.overrides</string>
			</dict>
			<key>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_ASR_ASSISTANT</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_ASR_ASSISTANT</string>
			</dict>
			<key>com.apple.voiceservices.notification.voice-update</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voiceservices.notification.voice-update</string>
			</dict>
			<key>com.apple.voicetrigger.EarlyDetect</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicetrigger.EarlyDetect</string>
			</dict>
			<key>com.apple.voicetrigger.PHSProfileModified</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicetrigger.PHSProfileModified</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.siri.xpc_activity.launchAtLogin</key>
			<dict>
				<key>Delay</key>
				<integer>0</integer>
				<key>GracePeriod</key>
				<integer>120</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.assistantd.fmip</key>
		<true/>
		<key>com.apple.assistant.analytics</key>
		<true/>
		<key>com.apple.assistant.analytics-observation</key>
		<true/>
		<key>com.apple.assistant.audio-session-assertion</key>
		<true/>
		<key>com.apple.assistant.client</key>
		<true/>
		<key>com.apple.assistant.crossdevice.election.service</key>
		<true/>
		<key>com.apple.assistant.dictation</key>
		<true/>
		<key>com.apple.assistant.flow_service</key>
		<true/>
		<key>com.apple.assistant.request-dispatcher.uibridge-service</key>
		<true/>
		<key>com.apple.assistant.security</key>
		<true/>
		<key>com.apple.assistant.settings</key>
		<true/>
		<key>com.apple.assistant.sync</key>
		<true/>
		<key>com.apple.assistant.uibridge-service</key>
		<true/>
		<key>com.apple.assistantd.aps</key>
		<true/>
		<key>com.apple.assistantd.managedstorage</key>
		<true/>
		<key>com.apple.private.alloy.siri.icloud-idswake</key>
		<true/>
		<key>com.apple.siri.analytics.assistant</key>
		<true/>
		<key>com.apple.siri.client_lite</key>
		<true/>
		<key>com.apple.siri.contextcollector</key>
		<true/>
		<key>com.apple.siri.external_request</key>
		<true/>
		<key>com.apple.siri.intelligence-tool-helper-service</key>
		<true/>
		<key>com.apple.siri.location</key>
		<true/>
		<key>com.apple.siri.orchestration.capabilities</key>
		<true/>
		<key>com.apple.siri.orchestration.intelligencesession</key>
		<true/>
		<key>com.apple.siri.orchestration.sirix</key>
		<true/>
		<key>com.apple.siri.pommes_search_xpc_service</key>
		<true/>
		<key>com.apple.siri.tether</key>
		<true/>
		<key>com.apple.siri.vocabularyupdates</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AssistantServices.framework/Versions/A/Support/assistantd</string>
	</array>
</dict>
</plist>
