# beta/macOS/System/Library/LaunchAgents/com.apple.biomesyncd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.biomesyncd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.cascade.xpc_event.setChange</key>
		<dict>
			<key>com.apple.biomesyncd.cascadeSetChange</key>
			<dict/>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.rapport.CompanionLinkDeviceAdded</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.rapport.CompanionLinkDeviceAdded</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>BiomeRapportWake</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.biomesyncd.rapport</string>
				<key>type</key>
				<string>server</string>
			</dict>
			<key>CascadeRapportWake</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.biomesyncd.cascade.rapport</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.biomesyncd.cascade.periodic-sync</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>ExpectedDuration</key>
				<real>300</real>
				<key>Interval</key>
				<real>3600</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<real>3600</real>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.biomesyncd.database-maintenance</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>ExpectedDuration</key>
				<real>300</real>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<real>3600</real>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<false/>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.biomesyncd.metrics-collection</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>ExpectedDuration</key>
				<real>300</real>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<real>3600</real>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.biomesyncd.periodic-sync</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>ExpectedDuration</key>
				<real>300</real>
				<key>Interval</key>
				<real>3600</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<real>3600</real>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.biomesyncd.post-install-immediate-sync</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>ExpectedDuration</key>
				<real>300</real>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.biome.ContextSync</key>
		<true/>
		<key>com.apple.biomesyncd.sync</key>
		<true/>
		<key>com.apple.private.alloy.contextsync-idswake</key>
		<true/>
		<key>com.apple.private.alloy.contextsync.local-idswake</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/biomesyncd</string>
	</array>
</dict>
</plist>
