# beta/macOS/System/Library/LaunchAgents/com.apple.FileProvider.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LIBDISPATCH_STRICT</key>
		<string>YES</string>
		<key>OBJC_DEBUG_MISSING_POOLS</key>
		<string>halt</string>
	</dict>
	<key>Label</key>
	<string>com.apple.FileProvider</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.FileProvider.fpfs.telemetry</key>
			<dict>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresUserInactivity</key>
				<false/>
			</dict>
		</dict>
		<key>com.apple.fsevents.matching</key>
		<dict>
			<key>com.apple.bird.fsevents</key>
			<dict>
				<key>Path</key>
				<string>~/Library/CloudStorage</string>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.FileProvider</key>
		<true/>
		<key>com.apple.FileProvider.cache-delete</key>
		<true/>
		<key>com.apple.FileProvider.sim2host</key>
		<true/>
		<key>com.apple.corespotlight.daemon.fileprovider</key>
		<true/>
		<key>com.apple.fileprovider.extension.push</key>
		<true/>
		<key>com.apple.fileprovider.extension.push.development</key>
		<true/>
		<key>com.apple.fileprovider.pushkit</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/FileProvider.framework/Support/fileproviderd</string>
</dict>
</plist>
