# beta/macOS/System/Library/LaunchAgents/com.apple.linkd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.linkd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>ApplicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>ApplicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.linkd.prune-transcript</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.linkd.update-registry</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.linkd.autoShortcut</key>
		<true/>
		<key>com.apple.linkd.cliSupport</key>
		<true/>
		<key>com.apple.linkd.constraints</key>
		<true/>
		<key>com.apple.linkd.extension</key>
		<true/>
		<key>com.apple.linkd.mediator</key>
		<true/>
		<key>com.apple.linkd.observationStatusRegistry</key>
		<true/>
		<key>com.apple.linkd.registry</key>
		<true/>
		<key>com.apple.linkd.suggestions</key>
		<true/>
		<key>com.apple.linkd.transcript</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/linkd</string>
	</array>
	<key>Umask</key>
	<integer>54</integer>
</dict>
</plist>
