# beta/macOS/System/Library/LaunchAgents/com.apple.CallHistorySyncHelper.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.CallHistorySyncHelper</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.callhistorysync.idslaunchnotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.callhistorysync.idslaunchnotification</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.CallHistorySyncHelper</key>
		<true/>
		<key>com.apple.CallHistorySyncHelper.aps</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CallHistory.framework/Support/CallHistorySyncHelper</string>
	</array>
</dict>
</plist>
