# beta/macOS/System/Library/LaunchAgents/com.apple.knowledge-agent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.apple.knowledge-agent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.coreduet.xpc.receiver</key>
		<dict>
			<key>appintent</key>
			<dict/>
		</dict>
		<key>com.apple.coreduetcontext.client_event_stream</key>
		<dict>
			<key>notification</key>
			<dict>
				<key>client-identifier</key>
				<string>com.apple.knowledge-agent</string>
			</dict>
		</dict>
		<key>com.apple.coreduetcontext.mdcs_events</key>
		<dict>
			<key>mdcs</key>
			<dict/>
		</dict>
		<key>com.apple.coreduetinteraction.interaction-deleted</key>
		<dict>
			<key>interaction-deleted</key>
			<dict/>
		</dict>
		<key>com.apple.coreduetinteraction.interaction-recorded</key>
		<dict>
			<key>interaction-recorded</key>
			<dict/>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>__ABDataBaseChangedByOtherProcessNotification</key>
			<dict>
				<key>Notification</key>
				<string>__ABDataBaseChangedByOtherProcessNotification</string>
			</dict>
			<key>com.apple.coreduet.client-needs-help.knowledge-agent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.coreduet.client-needs-help.knowledge-agent</string>
			</dict>
			<key>com.apple.rapport.CompanionLinkDeviceAdded</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.rapport.CompanionLinkDeviceAdded</string>
			</dict>
			<key>com.apple.siri.cloud.storage.deleted</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.cloud.storage.deleted</string>
			</dict>
			<key>com.apple.siri.cloud.synch.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.cloud.synch.changed</string>
			</dict>
			<key>com.apple.suggestions.settingsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.suggestions.settingsChanged</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>OnIncomingRequest</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.coreduet.sync</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.coreduet.knowledge.sync.push</key>
		<true/>
		<key>com.apple.coreduetd.knowledge.user</key>
		<true/>
		<key>com.apple.coreduetd.people.user</key>
		<true/>
		<key>com.apple.proactive.SuggestionRequest.cd_calendar_interaction_suggestions</key>
		<true/>
		<key>com.apple.proactive.SuggestionRequest.ps_facetime_fallback_interaction_model</key>
		<true/>
		<key>com.apple.proactive.SuggestionRequest.ps_facetime_interaction_model</key>
		<true/>
		<key>com.apple.proactive.SuggestionRequest.ps_unstructured_calendar_interaction_suggestions</key>
		<true/>
		<key>com.apple.proactive.SuggestionRequest.ps_unstructured_reminder_interaction_suggestions</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/knowledge-agent</string>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
