# beta/macOS/System/Library/LaunchAgents/com.apple.usernoted.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.usernoted</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.corespotlight.daemon.notifications</key>
		<true/>
		<key>com.apple.usernoted.client</key>
		<true/>
		<key>com.apple.usernoted.daemon_client</key>
		<true/>
		<key>com.apple.usernoted.dock</key>
		<true/>
		<key>com.apple.usernoted.events</key>
		<true/>
		<key>com.apple.usernoted.multi_client</key>
		<true/>
		<key>com.apple.usernoted.notificationcenter</key>
		<true/>
		<key>com.apple.usernoted.prefs</key>
		<true/>
		<key>com.apple.usernoted.push</key>
		<true/>
		<key>com.apple.usernoted.test</key>
		<true/>
		<key>com.apple.usernotifications.remotenotificationservice</key>
		<true/>
		<key>com.apple.usernotifications.usernotificationservice</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/usr/sbin/usernoted</string>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
