# beta/macOS/System/Library/LaunchAgents/com.apple.contactsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>HardResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>1024</integer>
	</dict>
	<key>Label</key>
	<string>com.apple.contactsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.private.contacts.accounts.validate-cache</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresUserInactivity</key>
				<true/>
				<key>ShouldWakeDevice</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>ABDatabaseChangedExternallyNotification</key>
			<dict>
				<key>Name</key>
				<string>ABDatabaseChangedExternallyNotification</string>
			</dict>
			<key>ABDatabaseChangedExternallyNotificationPriv</key>
			<dict>
				<key>Name</key>
				<string>ABDatabaseChangedExternallyNotificationPriv</string>
			</dict>
			<key>ABDistributedDatabaseChangedNotification</key>
			<dict>
				<key>Name</key>
				<string>ABDistributedDatabaseChangedNotification</string>
			</dict>
			<key>ABDistributedDatabaseChangedNotificationPriv</key>
			<dict>
				<key>Name</key>
				<string>ABDistributedDatabaseChangedNotificationPriv</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.contactsd.VerifyCoreSpotlightService</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.contactsd.accounts.cleanup</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<real>86400</real>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Background</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.AddressBook.ContactsAccountsService</key>
		<true/>
		<key>com.apple.askto.client.com.apple.contacts.commlimits</key>
		<true/>
		<key>com.apple.contacts.account-caching</key>
		<true/>
		<key>com.apple.contactsd.contact-provider</key>
		<true/>
		<key>com.apple.contactsd.persistence</key>
		<true/>
		<key>com.apple.contactsd.support</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/Contacts.framework/Support/contactsd</string>
	</array>
	<key>SoftResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>1024</integer>
	</dict>
</dict>
</plist>
