# beta/macOS/System/Library/LaunchAgents/com.apple.askpermissiond.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.askpermissiond</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>AppleLanguagePreferencesChangedNotification</key>
			<dict>
				<key>Name</key>
				<string>AppleLanguagePreferencesChangedNotification</string>
			</dict>
		</dict>
		<key>com.apple.usernotificationcenter.matching</key>
		<dict>
			<key>askpermissiond</key>
			<dict>
				<key>bundleid</key>
				<string>com.apple.askpermissiond</string>
				<key>events</key>
				<array>
					<string>didDismissAlert</string>
					<string>didActivateNotification</string>
					<string>didDeliverNotification</string>
					<string>didSnoozeAlert</string>
					<string>didRemoveDeliveredNotifications</string>
					<string>didExpireNotifications</string>
				</array>
				<key>system</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.askpermission</key>
		<true/>
		<key>com.apple.askpermissiond</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.askpermission.notifications</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/AskPermission.framework/Versions/A/Resources/askpermissiond</string>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
