# beta/macOS/System/Library/LaunchAgents/com.apple.GameController.gamecontrolleragentd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.GameController.gamecontrolleragentd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.gamecontroller.agent.launch</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.gamecontroller.agent.launch</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.GameController.gamecontrolleragentd</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/gamecontrolleragentd</string>
	</array>
</dict>
</plist>
