# beta/macOS/System/Library/LaunchAgents/com.apple.knowledgeconstructiond.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.knowledgeconstructiond</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.FullPipeline</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>ResourceIntensive</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.Pipeline.FastPass</key>
			<dict>
				<key>Dependencies</key>
				<array>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.cloudphotod.initialdownload.fastpass</string>
						<key>DependencyType</key>
						<string>DependencyTypeCompletion</string>
					</dict>
				</array>
				<key>FeatureCodes</key>
				<array>
					<integer>700</integer>
				</array>
				<key>NonRepeatingTask</key>
				<dict>
					<key>FastPassTask</key>
					<dict>
						<key>ProcessingTaskIdentifiers</key>
						<array>
							<string>com.apple.intelligenceplatform.IntelligencePlatformCore.FullPipeline</string>
						</array>
					</dict>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.AutonamingMessagesAggregationAccuracy</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>CPUIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>604800</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.AutonamingMessagesPipeline</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>CPUIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.intelligenceplatform.KnowledgeConstruction</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/IntelligencePlatformCore.framework/Versions/A/knowledgeconstructiond</string>
	</array>
</dict>
</plist>
