# beta/macOS/System/Library/LaunchAgents/com.apple.managedappdistributionagent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<dict>
		<key>#IfFeatureFlagDisabled</key>
		<string>ManagedAppDistribution/Service</string>
		<key>#Then</key>
		<true/>
	</dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.managedappdistributionagent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.managedappdistributionagent.ddm-repair</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>57600</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresUserInactivity</key>
				<false/>
			</dict>
			<key>com.apple.managedappdistributionagent.update-poll</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>57600</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.nsurlsessiond-events</key>
		<dict>
			<key>LaunchOnDemand</key>
			<dict>
				<key>ClientIdentifier</key>
				<string>com.apple.managedappdistributiond</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.managedappdistributionagent.xpc</key>
		<true/>
		<key>com.apple.managedappdistributionagent.xpc.legacy</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/Frameworks/ManagedAppDistribution.framework/Support/managedappdistributionagent</string>
</dict>
</plist>
