# beta/macOS/System/Library/LaunchAgents/com.apple.FeatureAccessAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.FeatureAccessAgent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>DeviceDidPairEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
			<key>DeviceNameChangeEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.lockdown.device_name_changed</string>
			</dict>
			<key>FirstUnlockEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>LanguageChangedEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.language.changed</string>
			</dict>
			<key>PhotosOptimizeStateChangedEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobileslideshow.PLNotificationKeepOriginalsChanged</string>
			</dict>
			<key>PhotosiCPLStateChangedEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobileslideshow.ICPLStateChanged</string>
			</dict>
			<key>SimulateVFSAlmostFullEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.cloud.quota.simulate.vfs.almostfull</string>
			</dict>
			<key>SimulateVFSNotFullEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.cloud.quota.simulate.vfs.notfull</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.featureaccess.reload</key>
			<dict>
				<key>AllowBattery</key>
				<string>1</string>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<string>1</string>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.ind.cloudfeatures</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CloudSubscriptionFeatures.framework/FeatureAccessAgent</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
