# beta/macOS/System/Library/LaunchAgents/com.apple.printtool.agent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PRINTTOOL_MODE</key>
		<string>agent</string>
	</dict>
	<key>Label</key>
	<string>com.apple.printtool.agent</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.printtool.agent</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/ApplicationServices.framework/Frameworks/PrintCore.framework/Versions/A/printtool</string>
	</array>
	<key>com.apple.private.security.restricted-application-groups</key>
	<array>
		<string>group.com.apple.printtool</string>
	</array>
	<key>com.apple.security.application-groups</key>
	<array>
		<string>group.com.apple.printtool</string>
	</array>
</dict>
</plist>
