# beta/macOS/System/Library/LaunchAgents/com.apple.menuextra.battery.helper.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.menuextra.battery.helper</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>LoginWindow</string>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.BMEHelper</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.battery-monitor</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>App</string>
	<key>Program</key>
	<string>/System/Library/CoreServices/Menu Extras/Battery.menu/Contents/Resources/bmehelper</string>
</dict>
</plist>
