# beta/macOS/System/Library/LaunchAgents/com.apple.diagnosticspushd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.diagnosticspushd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.diagnosticspushd.postinstall</key>
			<dict>
				<key>NonRepeatingTask</key>
				<dict>
					<key>ScheduleAfter</key>
					<integer>0</integer>
					<key>TrySchedulingBefore</key>
					<integer>300</integer>
				</dict>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>D</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.diagnostics</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/diagnosticspushd</string>
	</array>
</dict>
</plist>
