# beta/macOS/System/Library/LaunchAgents/com.apple.storeuid.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>ExitTimeOut</key>
	<integer>1</integer>
	<key>Label</key>
	<string>com.apple.storeuid</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>LoginWindow</string>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.storeagent.storekit</key>
		<true/>
		<key>com.apple.storeuid</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>_AdaptiveUtility</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/CommerceKit.framework/Versions/A/Resources/storeuid.app/Contents/MacOS/storeuid</string>
</dict>
</plist>
