# beta/macOS/System/Library/LaunchAgents/com.apple.replicatord.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LS_NO_IN_PROCESS_DATABASE_BUILD</key>
		<string>1</string>
		<key>MallocStackLoggingDirectory</key>
		<string>/private/var/tmp/com.apple.replicatord</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>Crashed</key>
		<true/>
	</dict>
	<key>Label</key>
	<string>com.apple.replicatord</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>keybag-first-unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>state-replicator-application-service</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.ApplicationService.chrono</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.private.alloy.widgets-idswake</key>
		<true/>
		<key>com.apple.replicatorservices</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/ReplicatorCore.framework/Support/replicatord</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>ThrottleInterval</key>
	<integer>3</integer>
</dict>
</plist>
