# beta/macOS/System/Library/LaunchAgents/com.apple.asktod.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.asktod</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.biome.compute.publisher.user</key>
		<dict>
			<key>com.apple.asktod.screenTimeRequest</key>
			<dict>
				<key>client</key>
				<string>com.apple.asktod</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>ScreenTimeRequest</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.asktod.screenTimeRequest</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.asktod</key>
		<true/>
		<key>com.apple.private.alloy.askto-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.asktod</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/asktod</string>
	</array>
</dict>
</plist>
