# beta/macOS/System/Library/LaunchAgents/com.apple.networkserviceproxy-osx.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.networkserviceproxy</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>ApplicationDatabaseUpdated</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.DatabaseSeedingComplete</string>
			</dict>
			<key>ApplicationInstalled</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>ApplicationUninstalled</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.networkserviceproxy</key>
		<dict>
			<key>SendMessage</key>
			<dict>
				<key>MessageID</key>
				<integer>0</integer>
			</dict>
		</dict>
		<key>com.apple.networkserviceproxy.file-descriptor-maintainer</key>
		<dict>
			<key>new-handles</key>
			<dict>
				<key>handles-array</key>
				<array/>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.CloudSubscriptionFeature.Changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.CloudSubscriptionFeature.Changed</string>
			</dict>
			<key>com.apple.networkextension.app-paths-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.networkextension.app-paths-changed</string>
			</dict>
			<key>com.apple.networkserviceproxy.reset</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.networkserviceproxy.reset</string>
			</dict>
			<key>com.apple.trial.NamespaceUpdate.NETWORK_SERVICE_PROXY_CONFIG_UPDATE</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.NETWORK_SERVICE_PROXY_CONFIG_UPDATE</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.networkserviceproxy.metadata-refresh-12h</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>43200</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.networkserviceproxy.metadata-refresh-22h</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>79200</integer>
				<key>PowerNap</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.networkserviceproxy.publish-daily-stats</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.networkserviceproxy</key>
		<true/>
		<key>com.apple.networkserviceproxy.fetch-token</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/networkserviceproxy</string>
</dict>
</plist>
