# beta/macOS/System/Library/LaunchAgents/com.apple.chronod.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LS_NO_IN_PROCESS_DATABASE_BUILD</key>
		<string>1</string>
		<key>MallocStackLoggingDirectory</key>
		<string>/private/var/tmp/com.apple.chrono</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>Crashed</key>
		<true/>
	</dict>
	<key>Label</key>
	<string>com.apple.chronod</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>application-installed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>application-restriction-changes</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationStateChanged</string>
			</dict>
			<key>application-uninstalled</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.chronod.replicator.message</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.chronod.replicator.message</string>
			</dict>
			<key>com.apple.chronod.replicator.record</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.chronod.replicator.record</string>
			</dict>
			<key>com.apple.replicatord.devicesChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.replicatord.devicesChanged</string>
			</dict>
			<key>environment-accessibility-bold-text-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.cache.enhance.text.legibility</string>
			</dict>
			<key>environment-accessibility-differentiate-without-color-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.cache.differentiate.without.color</string>
			</dict>
			<key>environment-accessibility-increased-contrast-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.cache.darken.system.colors.enabled</string>
			</dict>
			<key>environment-accessibility-invert-colors-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.cache.invert.colors</string>
			</dict>
			<key>environment-accessibility-preferred-content-size-changed</key>
			<dict>
				<key>Notification</key>
				<string>ApplePreferredContentSizeCategoryChangedNotification</string>
			</dict>
			<key>environment-accessibility-reduce-motion-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.cache.reduce.motion</string>
			</dict>
			<key>environment-accessibility-reduce-transparency-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.cache.enhance.background.contrast</string>
			</dict>
			<key>environment-locale-language-changed</key>
			<dict>
				<key>Notification</key>
				<string>AppleLanguagePreferencesChangedNotification</string>
			</dict>
			<key>environment-locale-region-changed</key>
			<dict>
				<key>Notification</key>
				<string>AppleNumberPreferencesChangedNotification</string>
			</dict>
			<key>environment-locale-region-date-changed</key>
			<dict>
				<key>Notification</key>
				<string>AppleDatePreferencesChangedNotification</string>
			</dict>
			<key>environment-locale-region-time-changed</key>
			<dict>
				<key>Notification</key>
				<string>AppleTimePreferencesChangedNotification</string>
			</dict>
			<key>keybag-status-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.lock_status</string>
			</dict>
			<key>significant-time-changed</key>
			<dict>
				<key>Notification</key>
				<string>SignificantTimeChangeNotification</string>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.chrono.PushNotifications</key>
		<true/>
		<key>com.apple.aps.chrono.PushNotifications.dev</key>
		<true/>
		<key>com.apple.chrono.controlcenter</key>
		<true/>
		<key>com.apple.chrono.widgetcenterconnection</key>
		<true/>
		<key>com.apple.chronod.replicator</key>
		<true/>
		<key>com.apple.chronoservices</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/ChronoCore.framework/Support/chronod</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>ThrottleInterval</key>
	<integer>3</integer>
</dict>
</plist>
