# beta/macOS/System/Library/LaunchAgents/com.apple.wallpaper.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<dict>
		<key>AfterInitialDemand</key>
		<true/>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.apple.wallpaper.agent</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.usernotifications.delegate.com.apple.wallpaper.notifications.sonoma-first-run</key>
		<true/>
		<key>com.apple.wallpaper</key>
		<true/>
		<key>com.apple.wallpaper.CacheDelete</key>
		<true/>
		<key>com.apple.wallpaper.debug.service</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>Program</key>
	<string>/System/Library/CoreServices/WallpaperAgent.app/Contents/MacOS/WallpaperAgent</string>
	<key>RunAtLoad</key>
	<false/>
	<key>ThrottleInterval</key>
	<integer>1</integer>
</dict>
</plist>
