# beta/macOS/System/Library/LaunchAgents/com.apple.SpacesTouchBarAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.SpacesTouchBarAgent.app</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.SpacesTouchBarAgent.touchbar.showPopover</key>
			<dict>
				<key>Name</key>
				<string>com.apple.SpacesTouchBarAgent.showPopover</string>
			</dict>
		</dict>
		<key>com.apple.touchbar.matching</key>
		<dict>
			<key>com.apple.touchbar.matching</key>
			<dict>
				<key>PresentationMode</key>
				<array>
					<string>spaces</string>
					<string>spacesWithControlStrip</string>
				</array>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.SpacesTouchBarService</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>App</string>
	<key>Program</key>
	<string>/System/Library/CoreServices/SpacesTouchBarAgent.app/Contents/MacOS/SpacesTouchBarAgent</string>
</dict>
</plist>
