# beta/macOS/System/Library/LaunchAgents/com.apple.passd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.passd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.eos.support-host</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IONotificationType</key>
				<string>IOServiceFirstMatch</string>
				<key>IOProviderClass</key>
				<string>AppleEmbeddedOSSupportHost</string>
			</dict>
		</dict>
		<key>com.apple.nsurlsessiond-events</key>
		<dict>
			<key>LaunchOnDemand</key>
			<dict>
				<key>ClientIdentifier</key>
				<string>com.apple.passd</string>
			</dict>
		</dict>
		<key>com.apple.usernotificationcenter.matching</key>
		<dict>
			<key>com.apple.passd</key>
			<dict>
				<key>centerType</key>
				<integer>6</integer>
				<key>events</key>
				<array>
					<string>didDeliverNotification</string>
					<string>didActivateNotification</string>
					<string>didRemoveDeliveredNotifications</string>
				</array>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.passd.account</key>
		<true/>
		<key>com.apple.passd.aps</key>
		<true/>
		<key>com.apple.passd.assertions</key>
		<true/>
		<key>com.apple.passd.bulletins</key>
		<true/>
		<key>com.apple.passd.cloud-store</key>
		<true/>
		<key>com.apple.passd.device-registration</key>
		<true/>
		<key>com.apple.passd.in-app-payment</key>
		<true/>
		<key>com.apple.passd.library</key>
		<true/>
		<key>com.apple.passd.payment</key>
		<true/>
		<key>com.apple.passd.payment-continuity</key>
		<true/>
		<key>com.apple.passd.sharing-channel</key>
		<true/>
		<key>com.apple.passd.trusted-device-enrollment-info-provider</key>
		<true/>
		<key>com.apple.passd.usage</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.private.alloy.applepay-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.Passbook</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/PassKitCore.framework/passd</string>
	</array>
</dict>
</plist>
