# beta/macOS/System/Library/LaunchAgents/com.apple.mbfloagent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.mbfloagent</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mbfloagent</key>
		<true/>
	</dict>
	<key>MultipleInstances</key>
	<true/>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/Setup Assistant.app/Contents/Resources/mbfloagent</string>
	</array>
</dict>
</plist>
