# beta/macOS/System/Library/LaunchAgents/com.apple.rapportd-user.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.rapportd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.rapport.prefsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.rapport.prefsChanged</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.CompanionLink</key>
		<true/>
		<key>com.apple.PairingManager</key>
		<true/>
		<key>com.apple.RemoteDisplay</key>
		<true/>
		<key>com.apple.private.alloy.coreduet.sync-idswake</key>
		<true/>
		<key>com.apple.private.alloy.nearby-idswake</key>
		<true/>
		<key>com.apple.rapport</key>
		<true/>
		<key>com.apple.rapport.RPPairing</key>
		<true/>
		<key>com.apple.rapport.StatusUpdates</key>
		<true/>
		<key>com.apple.rapport.people</key>
		<true/>
		<key>com.apple.rapport.private-discovery</key>
		<true/>
		<key>com.apple.rapport.remote-text-input</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/usr/libexec/rapportd</string>
	<key>PublishesEvents</key>
	<string>com.apple.rapport.matching</string>
</dict>
</plist>
