# beta/macOS/System/Library/LaunchAgents/com.apple.AquaAppearanceHelper.agent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.AquaAppearanceHelper</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aquaappearancehelper</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/SkyLight.framework//Versions/A/Resources/AquaAppearanceHelper.app/Contents/MacOS/AquaAppearanceHelper</string>
	</array>
</dict>
</plist>
