# beta/macOS/System/Library/LaunchAgents/com.apple.GamePolicyAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.GamePolicyAgent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.gamepolicy.agent.launch</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.gamepolicy.agent.launch</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.GamePolicyAgent.daemon</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.GamePolicyNotifications</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>App</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/GamePolicyAgent</string>
	</array>
</dict>
</plist>
