# beta/macOS/System/Library/LaunchAgents/com.apple.CloudSettingsSyncAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.cloudsettingssyncagent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.fakebuddy.launchnew</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.fakebuddy.launchnew</string>
			</dict>
			<key>com.apple.fakebuddy.launchupdate</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.fakebuddy.launchupdate</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.appearance</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.appearance</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.controlcenter</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.controlcenter</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.datetime</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.datetime</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.desktop</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.desktop</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.displays</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.displays</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.gamecontroller</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.gamecontroller</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.general</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.general</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.homescreen</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.homescreen</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.international</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.international</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.keyboard</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.keyboard</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.languageregion</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.languageregion</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.mouse</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.mouse</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.screensaver</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.screensaver</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.sound</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.sound</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.cloudsettings.trackpad</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.cloudsettings.trackpad</string>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/CloudSettingsSyncAgent</string>
	</array>
</dict>
</plist>
