# beta/macOS/System/Library/LaunchAgents/com.apple.protectedcloudstorage.protectedcloudkeysyncing.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.protectedcloudstorage.protectedcloudkeysyncing</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>_CDPWalrusStateChangeDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>_CDPWalrusStateChangeDarwinNotification</string>
			</dict>
			<key>com.apple.ProtectedCloudStorage.mobileBackupStateChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ProtectedCloudStorage.mobileBackupStateChange</string>
			</dict>
			<key>com.apple.ProtectedCloudStorage.rollBackupDisabled</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ProtectedCloudStorage.rollBackupDisabled</string>
			</dict>
			<key>com.apple.ProtectedCloudStorage.rollIfAged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ProtectedCloudStorage.rollIfAged</string>
			</dict>
			<key>com.apple.ProtectedCloudStorage.rollNow</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ProtectedCloudStorage.rollNow</string>
			</dict>
			<key>com.apple.ProtectedCloudStorage.test.mobileBackupStateChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ProtectedCloudStorage.test.mobileBackupStateChange</string>
			</dict>
			<key>com.apple.ProtectedCloudStorage.updatedKeys</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ProtectedCloudStorage.updatedKeys</string>
			</dict>
			<key>com.apple.security.secureobjectsync.circlechanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.security.secureobjectsync.circlechanged</string>
			</dict>
			<key>com.apple.security.view-change.PCS</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.security.view-change.PCS</string>
			</dict>
			<key>com.apple.security.view-ready.Manatee</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.security.view-ready.Manatee</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.ProtectedCloudStorage.fortnightly</key>
			<dict>
				<key>Interval</key>
				<integer>1209600</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
			<key>com.apple.ProtectedCloudStorage.scheduledKeyroll</key>
			<dict>
				<key>Delay</key>
				<integer>60</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>3600</integer>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.private.alloy.pcskey.sync-idswake</key>
		<true/>
		<key>com.apple.protectedcloudstorage.protectedcloudkeysyncing</key>
		<true/>
		<key>com.apple.protectedcloudstorage.protectedcloudkeysyncing.pairedsync</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/ProtectedCloudStorage.framework/Helpers/ProtectedCloudKeySyncing</string>
</dict>
</plist>
