# beta/macOS/System/Library/LaunchAgents/com.apple.watchlistd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MallocMaxMagazines</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.watchlistd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.AppleMediaServices.deviceOffersChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AppleMediaServices.deviceOffersChanged</string>
			</dict>
			<key>com.apple.VideoSubscriberAccount.DidRegisterSubscription</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.VideoSubscriberAccount.DidRegisterSubscription</string>
			</dict>
			<key>com.apple.VideosUI.PlayHistoryUpdatedNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.VideosUI.PlayHistoryUpdatedNotification</string>
			</dict>
			<key>com.apple.VideosUI.StoreAcquisitionCrossProcessNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.VideosUI.StoreAcquisitionCrossProcessNotification</string>
			</dict>
			<key>com.apple.VideosUI.UpNextRequestDidFinishNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.VideosUI.UpNextRequestDidFinishNotification</string>
			</dict>
			<key>com.apple.tv.TVWidgetExtension.Register</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.tv.TVWidgetExtension.Register</string>
			</dict>
			<key>com.apple.tv.appRemoved</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.tv.appRemoved</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.watchlistd</key>
		<true/>
		<key>com.apple.aps.watchlistd.dev</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.tv</key>
		<true/>
		<key>com.apple.watchlistd.aps.demo</key>
		<true/>
		<key>com.apple.watchlistd.xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/WatchListKit.framework/Support/watchlistd</string>
	</array>
	<key>com.apple.private.tcc.allow</key>
	<array>
		<string>kTCCServiceMediaLibrary</string>
	</array>
	<key>com.apple.security.exception.mach-lookup.global-name</key>
	<array>
		<string>com.apple.usernotifications.usernotificationservice</string>
		<string>com.apple.usernotifications.listener</string>
	</array>
</dict>
</plist>
