# beta/macOS/System/Library/LaunchAgents/com.apple.usbnotificationagent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.usbnotificationagent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.sop.displayPort</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyExistsMatch</key>
				<array>
					<string>DP Port Hint</string>
				</array>
				<key>IOProviderClass</key>
				<string>IOPortTransportComponentCCUSBPDSOP</string>
			</dict>
			<key>com.apple.usb.billboardFailedDevice</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyExistsMatch</key>
				<array>
					<string>UsbBillboardAltModeFailed</string>
				</array>
				<key>IOProviderClass</key>
				<string>AppleUSBHostBillboardDevice</string>
			</dict>
			<key>com.apple.usb.deviceLinkSpeedFailed</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyExistsMatch</key>
				<array>
					<string>Usb3LinkRequired</string>
				</array>
				<key>IOProviderClass</key>
				<string>IOUSBHostDevice</string>
			</dict>
			<key>com.apple.usb.deviceLinkSpeedNonoptimal</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyExistsMatch</key>
				<array>
					<string>Usb3LinkPreferred</string>
				</array>
				<key>IOProviderClass</key>
				<string>IOUSBHostDevice</string>
			</dict>
			<key>com.apple.usb.powerFailedDevice</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyExistsMatch</key>
				<array>
					<string>kUSBFailedRequestedPower</string>
				</array>
				<key>IOProviderClass</key>
				<string>IOUSBHostDevice</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.usernotifications.delegate.com.apple.usbnotificationagent</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/usbnotificationagent</string>
</dict>
</plist>
