# beta/macOS/System/Library/LaunchAgents/com.apple.AppSSOAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.AppSSOAgent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.AppSSO.configLoaded</key>
			<dict>
				<key>Name</key>
				<string>com.apple.AppSSO.configLoaded</string>
			</dict>
			<key>com.apple.KerberosExtension.configRemoved</key>
			<dict>
				<key>Name</key>
				<string>com.apple.KerberosExtension.configRemoved</string>
			</dict>
			<key>com.apple.LaunchServices.pluginsregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.pluginsregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.Kerberos.cache.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.Kerberos.cache.changed</string>
			</dict>
			<key>com.apple.PlatformSSO.unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.PlatformSSO.unlock</string>
			</dict>
			<key>com.apple.network_change</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.config.network_change.dns</string>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.AppSSO.service-xpc</key>
		<dict/>
		<key>com.apple.PlatformSSO.service-login-manager-xpc</key>
		<dict/>
		<key>com.apple.PlatformSSO.service-xpc</key>
		<dict/>
		<key>com.apple.PlatformSSO.settings.service-xpc</key>
		<dict/>
		<key>com.apple.usernotifications.delegate.com.apple.PlatformSSO.notifications</key>
		<dict/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AppSSO.framework/Support/AppSSOAgent.app/Contents/MacOS/AppSSOAgent</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
