# beta/macOS/System/Library/LaunchAgents/com.apple.peopled.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>ExitTimeOut</key>
	<integer>1</integer>
	<key>Label</key>
	<string>com.apple.peopled</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.biome.compute.publisher.user</key>
		<dict>
			<key>com.apple.people.askToBuyRequest</key>
			<dict>
				<key>client</key>
				<string>com.apple.peopled</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>AskToBuy</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.people.askToBuyRequest</string>
			</dict>
			<key>com.apple.people.childState</key>
			<dict>
				<key>client</key>
				<string>com.apple.peopled</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>Family.ScreenTime.ChildState</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.people.childState</string>
			</dict>
			<key>com.apple.people.dndAvailability</key>
			<dict>
				<key>client</key>
				<string>com.apple.peopled</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>UserStatusChange</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.people.dndAvailability</string>
			</dict>
			<key>com.apple.people.location</key>
			<dict>
				<key>client</key>
				<string>com.apple.peopled</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>FindMyLocationChange</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.people.location</string>
			</dict>
			<key>com.apple.people.missedCall</key>
			<dict>
				<key>client</key>
				<string>com.apple.peopled</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>AppIntent</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.people.missedCall</string>
			</dict>
			<key>com.apple.people.readMessage</key>
			<dict>
				<key>client</key>
				<string>com.apple.peopled</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>ReadMessage</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.people.readMessage</string>
			</dict>
			<key>com.apple.people.screenTimeRequest</key>
			<dict>
				<key>client</key>
				<string>com.apple.peopled</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>ScreenTimeRequest</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.people.screenTimeRequest</string>
			</dict>
			<key>com.apple.people.unreadMessage</key>
			<dict>
				<key>client</key>
				<string>com.apple.peopled</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>AppIntent</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.people.unreadMessage</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>_CalDatabaseChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>_CalDatabaseChangedNotification</string>
			</dict>
			<key>__ABDataBaseChangedByOtherProcessNotification</key>
			<dict>
				<key>Notification</key>
				<string>__ABDataBaseChangedByOtherProcessNotification</string>
			</dict>
			<key>com.apple.accessibility.cache.prefers.horizontal.text</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.cache.prefers.horizontal.text</string>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.corespotlight.daemon.people</key>
		<true/>
		<key>com.apple.people.agent</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/People.framework/peopled</string>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
