# beta/macOS/System/Library/LaunchAgents/com.apple.ap.adprivacyd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.ap.adprivacyd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CKAccountChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>CKAccountChangedNotification</string>
			</dict>
			<key>CPInstalledProfilesDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>CPInstalledProfilesDidChangeNotification</string>
			</dict>
			<key>com.apple.ap.adprivacyd.canceltasks</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ap.adprivacyd.canceltasks</string>
			</dict>
			<key>com.apple.ap.adprivacyd.iTunesActiveAccountDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ap.adprivacyd.iTunesActiveAccountDidChangeNotification</string>
			</dict>
			<key>com.apple.ap.adprivacyd.iTunesActiveStorefrontDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ap.adprivacyd.iTunesActiveStorefrontDidChangeNotification</string>
			</dict>
			<key>com.apple.ap.adprivacyd.launch</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ap.adprivacyd.launch</string>
			</dict>
			<key>com.apple.ap.adprivacyd.reconcile</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ap.adprivacyd.reconcile</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.ap.adprivacyd.reconcile</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>300</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequireScreenSleep</key>
				<false/>
				<key>RequiresBuddyComplete</key>
				<true/>
			</dict>
			<key>com.apple.ap.adprivacyd.rotateDeviceNewsPlusSubscriberID</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>300</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<false/>
				<key>RequireScreenSleep</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.AdSheetPhone.push</key>
		<true/>
		<key>com.apple.ap.adprivacyd.attribution</key>
		<true/>
		<key>com.apple.ap.adprivacyd.idmanager</key>
		<true/>
		<key>com.apple.ap.adprivacyd.opt-out</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/adprivacyd</string>
	<key>StartInterval</key>
	<integer>259200</integer>
</dict>
</plist>
