# beta/macOS/System/Library/LaunchAgents/com.apple.amp.mediasharingd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>Crashed</key>
		<true/>
	</dict>
	<key>Label</key>
	<string>com.apple.amp.mediasharingd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.amp.mediasharingd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AMPSharing.framework/Versions/A/Support/mediasharingd</string>
		<string>--launchd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
</dict>
</plist>
