# beta/macOS/System/Library/LaunchAgents/com.apple.cloudpaird.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>Crashed</key>
		<true/>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.apple.BTServer.cloudpairing</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bluetooth.discovery</key>
		<dict>
			<key>NearbyAudioAccessory</key>
			<dict>
				<key>bleRSSIThresholdHint</key>
				<integer>-50</integer>
				<key>discoveryFlags</key>
				<array>
					<string>WxStatus</string>
					<string>AccessoryStatus</string>
				</array>
				<key>label</key>
				<string>NearbyAudioAccessory</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CKIdentityUpdateNotification</key>
			<dict>
				<key>Notification</key>
				<string>CKIdentityUpdateNotification</string>
			</dict>
			<key>com.apple.coreaudio.BorealisToggled</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.coreaudio.BorealisToggled</string>
			</dict>
			<key>com.apple.icloudpairing.idslaunchnotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloudpairing.idslaunchnotification</string>
			</dict>
			<key>com.apple.system.hostname</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.hostname</string>
			</dict>
			<key>com.apple.xpc.activity</key>
			<dict>
				<key>com.apple.bluetooth.services.cloud.retryFetch</key>
				<dict>
					<key>AllowBattery</key>
					<true/>
					<key>PowerNap</key>
					<true/>
					<key>Priority</key>
					<string>Utility</string>
					<key>Repeating</key>
					<false/>
					<key>RequireNetworkConnectivity</key>
					<true/>
					<key>RequiresClassC</key>
					<true/>
				</dict>
				<key>com.apple.bluetooth.services.devicestore.deferredActions</key>
				<dict>
					<key>AllowBattery</key>
					<true/>
					<key>PowerNap</key>
					<true/>
					<key>Priority</key>
					<string>Utility</string>
					<key>Repeating</key>
					<false/>
					<key>RequireNetworkConnectivity</key>
					<true/>
					<key>RequiresClassC</key>
					<true/>
				</dict>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.AudioAccessoryServices</key>
		<true/>
		<key>com.apple.BTServer.cloudpairing</key>
		<true/>
		<key>com.apple.BluetoothCloudServices</key>
		<true/>
		<key>com.apple.BluetoothServices</key>
		<true/>
		<key>com.apple.HearingModeService</key>
		<true/>
		<key>com.apple.aps.cloudpaird</key>
		<true/>
		<key>com.apple.private.alloy.icloudpairing-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.AudioAccessoryUserNotifications</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.BTUserNotifications</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.MuteControlUserNotifications</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.SleepDetectionUserNotification</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/CoreServices/audioaccessoryd</string>
	<key>PublishesEvents</key>
	<array>
		<string>com.apple.AudioAccessoryDiscovery</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
