# beta/macOS/System/Library/LaunchAgents/com.apple.timezoneupdates.tznotify.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.timezoneupdates.tznotify</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.timezoneupdates.installednotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.timezoneupdates.installednotification</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.usernotifications.delegate.com.apple.TimeZoneUpdates.tznotify</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/tznotify</string>
</dict>
</plist>
