# beta/macOS/System/Library/LaunchAgents/com.apple.appstoreagent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.appstoreagent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.appstored.ASDItunesAccountDidChangeNotification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.appstored.ASDItunesAccountDidChangeNotification</string>
			</dict>
			<key>com.apple.appstored.ASDRepairBundleID</key>
			<dict>
				<key>Name</key>
				<string>com.apple.appstored.ASDRepairBundleID</string>
			</dict>
			<key>com.apple.appstored.ASDSandboxAccountDidChangeNotification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.appstored.ASDSandboxAccountDidChangeNotification</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.StoreServices.StorefrontChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.StoreServices.StorefrontChanged</string>
			</dict>
			<key>com.apple.crashreporter.third_party_submission_preference_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.crashreporter.third_party_submission_preference_changed</string>
			</dict>
			<key>com.apple.fairplayd.resync-fpkeybag</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.fairplayd.resync-fpkeybag</string>
			</dict>
		</dict>
		<key>com.apple.usernotificationcenter.matching</key>
		<dict>
			<key>appstoreagent</key>
			<dict>
				<key>bundleid</key>
				<string>com.apple.appstoreagent</string>
				<key>events</key>
				<array>
					<string>didDismissAlert</string>
					<string>didActivateNotification</string>
					<string>didDeliverNotification</string>
					<string>didSnoozeAlert</string>
					<string>didRemoveDeliveredNotifications</string>
					<string>didExpireNotifications</string>
				</array>
				<key>system</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.appstoreagent.scheduleTasks</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>120</integer>
				<key>GracePeriod</key>
				<integer>60</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<false/>
				<key>RequiresBuddyComplete</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.appstoreagent.CacheDelete</key>
		<true/>
		<key>com.apple.appstoreagent.legacy</key>
		<true/>
		<key>com.apple.appstoreagent.xpc</key>
		<true/>
		<key>com.apple.appstored.xpc.storequeue</key>
		<true/>
		<key>com.apple.aps.appstoreagent</key>
		<true/>
		<key>com.apple.ondemandd.CacheDelete</key>
		<true/>
		<key>com.apple.ondemandd.appinstall</key>
		<true/>
		<key>com.apple.ondemandd.appreview</key>
		<true/>
		<key>com.apple.ondemandd.client</key>
		<true/>
		<key>com.apple.ondemandd.devtools</key>
		<true/>
		<key>com.apple.ondemandd.launchservices</key>
		<true/>
		<key>com.apple.ondemandd.odrctl</key>
		<true/>
		<key>com.apple.storekitservice</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.AppStore</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/AppStoreDaemon.framework/Support/appstoreagent</string>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
