# beta/macOS/System/Library/LaunchAgents/com.apple.accessibility.heard.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.apple.accessibility.heard</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.accessibility.heard</key>
		<true/>
		<key>com.apple.aps.heard</key>
		<true/>
		<key>com.apple.private.alloy.hearing-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.SoundDetectionNotifications</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/HearingCore.framework/heard</string>
	<key>RunAtLoad</key>
	<false/>
	<key>ThrottleInterval</key>
	<integer>1</integer>
</dict>
</plist>
