# beta/macOS/System/Library/LaunchAgents/com.apple.symptomsd.distributed-agent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.symptomsd.distributed-agent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.symptoms.distributed.launch</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.symptoms.distributed.launch</string>
			</dict>
		</dict>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/symptomsd-distributed</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/symptomsd-distributed</string>
		<string>--agent</string>
	</array>
</dict>
</plist>
