# beta/macOS/System/Library/LaunchAgents/com.apple.DataDetectorsLocalSources.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>CFBundleIdentifier</key>
	<string>com.apple.DataDetectorsLocalSources</string>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.DataDetectorsLocalSources</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.DataDetectorsLocalSource</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<real>604800</real>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>POSIXSpawnType</key>
	<string>Background</string>
	<key>Program</key>
	<string>/usr/libexec/DataDetectorsLocalSources</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/DataDetectorsLocalSources</string>
	</array>
</dict>
</plist>
