# beta/macOS/System/Library/LaunchAgents/com.apple.EscrowSecurityAlert.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.EscrowSecurityAlert</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.EscrowSecurityAlert.record</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.EscrowSecurityAlert.record</string>
			</dict>
			<key>com.apple.EscrowSecurityAlert.reset</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.EscrowSecurityAlert.reset</string>
			</dict>
			<key>com.apple.EscrowSecurityAlert.server</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.EscrowSecurityAlert.server</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.EscrowSecurityAlert.daily</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/EscrowSecurityAlert.app/Contents/MacOS/EscrowSecurityAlert</string>
	</array>
</dict>
</plist>
