# beta/macOS/System/Library/LaunchAgents/com.apple.storekitagent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.storekitagent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.appstored.ASDItunesAccountDidChangeNotification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.appstored.ASDItunesAccountDidChangeNotification</string>
			</dict>
			<key>com.apple.appstored.ASDSandboxAccountDidChangeNotification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.appstored.ASDSandboxAccountDidChangeNotification</string>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.storekitservice</key>
		<true/>
		<key>com.apple.storekit.configuration.xpc</key>
		<true/>
		<key>com.apple.storekitagent</key>
		<true/>
		<key>com.apple.storekitservice.sim2host</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/Frameworks/StoreKit.framework/Support/storekitagent</string>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
