# beta/macOS/System/Library/LaunchAgents/com.apple.icloud.searchpartyuseragent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.icloud.searchpartyuseragent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CKIdentityUpdateNotification</key>
			<dict>
				<key>Notification</key>
				<string>CKIdentityUpdateNotification</string>
			</dict>
			<key>FMDSecureLocationsPushCommandReceived</key>
			<dict>
				<key>Notification</key>
				<string>FMDSecureLocationsReceivedCommand</string>
			</dict>
			<key>FMDSecureLocationsPushLocationReceived</key>
			<dict>
				<key>Notification</key>
				<string>FMDSecureLocationsReceivedPayload</string>
			</dict>
			<key>FMLFollowersChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>FMLFollowersChangedNotification</string>
			</dict>
			<key>com.apple.AOSNotification.FMIPStateDidChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AOSNotification.FMIPStateDidChange</string>
			</dict>
			<key>com.apple.FindMyDevice.LocatableStateReported</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.FindMyDevice.LocatableStateReported</string>
			</dict>
			<key>com.apple.family.family_updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.family.family_updated</string>
			</dict>
			<key>com.apple.icloud.searchpartyd.fetchServerConfiguration</key>
			<array>
				<string>com.apple.icloud.searchpartyd.fetchServerConfiguration</string>
			</array>
			<key>com.apple.icloud.spowner.evaluatelocalbeaconing</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.spowner.evaluatelocalbeaconing</string>
			</dict>
			<key>com.apple.locationd.authorization</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.authorization</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.icloud.searchpartyuseragent.post-install</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.searchpartyuseragent</key>
		<true/>
		<key>com.apple.icloud.searchparty.delegatedlocationpublish.items</key>
		<true/>
		<key>com.apple.icloud.searchparty.locationfetch.items</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.aps</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.beaconsharingservice</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.securelocations</key>
		<true/>
		<key>com.apple.icloud.searchpartyuseragent.beaconmanager</key>
		<true/>
		<key>com.apple.icloud.searchpartyuseragent.beaconmanager.simplebeacon</key>
		<true/>
		<key>com.apple.icloud.searchpartyuseragent.ownersession</key>
		<true/>
		<key>com.apple.icloud.searchpartyuseragent.pairingmanager</key>
		<true/>
		<key>com.apple.private.alloy.findmy.itemsharing-crossaccount-idswake</key>
		<true/>
		<key>com.apple.searchparty.btfindingsession</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/searchpartyuseragent</string>
	<key>Umask</key>
	<integer>23</integer>
</dict>
</plist>
