# beta/macOS/System/Library/LaunchAgents/com.apple.proactived.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.proactived</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.biome.compute.publisher</key>
		<dict>
			<key>FocusModes.GameController</key>
			<dict>
				<key>client</key>
				<string>com.apple.proactived</string>
				<key>graph</key>
				<dict>
					<key>upstream</key>
					<dict>
						<key>name</key>
						<string>streamPublisher</string>
						<key>streamIdentifier</key>
						<string>GameController</string>
						<key>version</key>
						<integer>1</integer>
					</dict>
				</dict>
				<key>identifier</key>
				<string>FocusModes.GameController</string>
			</dict>
			<key>FocusModes.Routine</key>
			<dict>
				<key>client</key>
				<string>com.apple.proactived</string>
				<key>graph</key>
				<dict>
					<key>upstream</key>
					<dict>
						<key>name</key>
						<string>streamPublisher</string>
						<key>streamIdentifier</key>
						<string>SemanticLocation</string>
						<key>version</key>
						<integer>1</integer>
					</dict>
				</dict>
				<key>identifier</key>
				<string>FocusModes.Routine</string>
			</dict>
			<key>FocusModes.ScreenShare</key>
			<dict>
				<key>client</key>
				<string>com.apple.proactived</string>
				<key>graph</key>
				<dict>
					<key>upstream</key>
					<dict>
						<key>name</key>
						<string>streamPublisher</string>
						<key>streamIdentifier</key>
						<string>ScreenSharing</string>
						<key>version</key>
						<integer>1</integer>
					</dict>
				</dict>
				<key>identifier</key>
				<string>FocusModes.ScreenShare</string>
			</dict>
			<key>FocusModes.WorkMode</key>
			<dict>
				<key>client</key>
				<string>com.apple.proactived</string>
				<key>graph</key>
				<dict>
					<key>upstream</key>
					<dict>
						<key>name</key>
						<string>streamPublisher</string>
						<key>streamIdentifier</key>
						<string>AppLaunch</string>
						<key>version</key>
						<integer>1</integer>
					</dict>
				</dict>
				<key>identifier</key>
				<string>FocusModes.WorkMode</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>_CalDatabaseChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>_CalDatabaseChangedNotification</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict/>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.proactive.ProactiveSuggestionClientModel.xpc</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/proactived</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/proactived</string>
	</array>
</dict>
</plist>
