# beta/Library/NanoLaunchDaemons/com.apple.companioncamerad.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.companioncamerad</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.companion.camera</key>
		<true/>
		<key>com.apple.private.alloy.camera.proxy-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CompanionCamera.framework/Support/companioncamerad</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>WorkingDirectory</key>
	<string>/tmp</string>
</dict>
</plist>
