# beta/Library/NanoLaunchDaemons/com.apple.nanophone.server.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.nanophone</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CPActiveCountryCodeChanged</key>
			<dict>
				<key>Notification</key>
				<string>CPActiveCountryCodeChanged</string>
			</dict>
			<key>CPHomeCountryCodeChanged</key>
			<dict>
				<key>Notification</key>
				<string>CPHomeCountryCodeChanged</string>
			</dict>
			<key>CPNetworkCountryCodeChanged</key>
			<dict>
				<key>Notification</key>
				<string>CPNetworkCountryCodeChanged</string>
			</dict>
			<key>TUReplyWithMessageStoreReplyArray</key>
			<dict>
				<key>Notification</key>
				<string>TUReplyWithMessageStoreReplyArray</string>
			</dict>
			<key>com.apple.FaceTimeMessageStore.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.FaceTimeMessageStore.changed</string>
			</dict>
			<key>com.apple.TelephonyUtilities.TUReplyWithMessageStoreChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TelephonyUtilities.TUReplyWithMessageStoreChanged</string>
			</dict>
			<key>com.apple.nanophone.vm.fullsync</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanophone.vm.fullsync</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.nanophone.pairedsync</key>
		<true/>
		<key>com.apple.private.alloy.voicemailsync-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/tursd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
