# beta/Library/NanoLaunchDaemonsAltAccount/com.apple.NPKCompanionAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>NSRunningFromLaunchd</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.NPKCompanionAgent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>BLTWillSendBulletinToGizmoNotification</key>
			<dict>
				<key>Name</key>
				<string>BLTWillSendBulletinToGizmoNotification</string>
			</dict>
			<key>PDSharedPeerPaymentWebServiceDidChangeNotification</key>
			<dict>
				<key>Name</key>
				<string>PDSharedPeerPaymentWebServiceDidChangeNotification</string>
			</dict>
			<key>PKPassLibraryCatalogDidChange</key>
			<dict>
				<key>Name</key>
				<string>PKPassLibraryCatalogDidChangeNotification</string>
			</dict>
			<key>PKPassLibraryDidAddPass</key>
			<dict>
				<key>Name</key>
				<string>PKPassLibraryDidAddPassNotification</string>
			</dict>
			<key>PKPassLibraryDidRemovePass</key>
			<dict>
				<key>Name</key>
				<string>PKPassLibraryDidRemovePassNotification</string>
			</dict>
			<key>PKPassLibraryDidUpdatePass</key>
			<dict>
				<key>Name</key>
				<string>PKPassLibraryDidUpdatePassNotification</string>
			</dict>
			<key>PKPassLibraryRelevantPassesDidChange</key>
			<dict>
				<key>Name</key>
				<string>PKPassLibraryRelevantPassesDidChangeNotification</string>
			</dict>
			<key>PKPassSettingsDidChangeNotification</key>
			<dict>
				<key>Name</key>
				<string>PKPassSettingsDidChangeNotification</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.AppSupport.loggingDefaultsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AppSupport.loggingDefaultsChanged</string>
			</dict>
			<key>com.apple.managedconfiguration.defaultsdidchange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.defaultsdidchange</string>
			</dict>
			<key>com.apple.passbook.globalprefschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.passbook.globalprefschanged</string>
			</dict>
			<key>com.apple.passbook.idslaunchnotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.passbook.idslaunchnotification</string>
			</dict>
			<key>com.apple.passkit.paymentoptionsdefaultschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.passkit.paymentoptionsdefaultschanged</string>
			</dict>
			<key>com.apple.stockholm.preferenceschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.stockholm.preferenceschanged</string>
			</dict>
		</dict>
		<key>com.apple.nsurlsessiond-events</key>
		<dict>
			<key>LaunchOnDemand</key>
			<dict>
				<key>ClientIdentifier</key>
				<string>com.apple.NPKCompanionAgent</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.NPKCompanionAgent.Server</key>
		<true/>
		<key>com.apple.NPKCompanionAgent.bulletins</key>
		<true/>
		<key>com.apple.NPKCompanionAgent.library</key>
		<true/>
		<key>com.apple.NPKCompanionAgent.pairedsync</key>
		<true/>
		<key>com.apple.NanoPassbook.IDVRemoteDeviceService.session.server</key>
		<true/>
		<key>com.apple.NanoPassbook.NPKCompanionViewService.connection.server</key>
		<true/>
		<key>com.apple.private.alloy.passbook.general-idswake</key>
		<true/>
		<key>com.apple.private.alloy.passbook.maintenance-idswake</key>
		<true/>
		<key>com.apple.private.alloy.passbook.relevancy-idswake</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/NanoPassKit.framework/NPKCompanionAgent</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
