# beta/macOS/System/Library/LaunchDaemons/com.apple.backgroundtaskmanagementd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.backgroundtaskmanagementd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.backgroundtaskmanagement</key>
		<true/>
		<key>com.apple.backgroundtaskmanagement.registrar</key>
		<true/>
		<key>com.apple.backgroundtaskmanagement.responses</key>
		<true/>
		<key>com.apple.backgroundtaskmanagement.sfl</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/BackgroundTaskManagement.framework/Resources/backgroundtaskmanagementd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/BackgroundTaskManagement.framework/Resources/backgroundtaskmanagementd</string>
		<string>-daemon</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
