# beta/macOS/System/Library/LaunchDaemons/com.apple.symptomsd-diag.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.symptomsd-diag</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.diagnosticpipeline.cloudkit_activity</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>NetworkParameters</key>
				<dict>
					<key>prohibit_constrained</key>
					<true/>
				</dict>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>NetworkTransferSize</key>
				<real>10485760</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.diagnosticpipeline.cloudkit_expedited_upload</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<real>14400</real>
				<key>NetworkParameters</key>
				<dict>
					<key>prohibit_constrained</key>
					<true/>
				</dict>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>NetworkTransferSize</key>
				<real>10485760</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.diagnosticpipeline.cloudkit_report_stats</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>NetworkParameters</key>
				<dict>
					<key>prohibit_constrained</key>
					<true/>
				</dict>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>NetworkTransferSize</key>
				<integer>1048576</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.diagnosticpipeline.tasking_invalidation_activity</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>7200</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.diagnosticpipeline.tasking</key>
		<true/>
		<key>com.apple.diagnosticpipeline.service</key>
		<true/>
		<key>com.apple.diagnosticpipeline.tasking.service</key>
		<true/>
		<key>com.apple.symptom_diagnostics</key>
		<true/>
		<key>com.apple.symptomsd-diag.CacheDelete</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/symptomsd-diag</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/symptomsd-diag</string>
	</array>
	<key>PublishesEvents</key>
	<string>com.apple.diagnosticpipeline.taskingevents</string>
</dict>
</plist>
