# beta/macOS/System/Library/LaunchDaemons/com.apple.bluetoothd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>OS_ACTIVITY_ENABLE_DYNAMIC</key>
		<string>1</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.apple.bluetoothd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.mobile.lockdown.activation_state</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.lockdown.activation_state</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.BTAudioHALPlugin.xpc</key>
		<true/>
		<key>com.apple.ExposureNotification</key>
		<true/>
		<key>com.apple.ULLAConcierge.xpc.diagnostics</key>
		<true/>
		<key>com.apple.bluetooth.BTPacketLogger</key>
		<true/>
		<key>com.apple.bluetooth.BTPacketLogger.ios</key>
		<true/>
		<key>com.apple.bluetooth.cloudkit.xpc</key>
		<true/>
		<key>com.apple.bluetooth.xpc</key>
		<true/>
		<key>com.apple.server.bluetooth</key>
		<true/>
		<key>com.apple.server.bluetooth.bbprovider.xpc</key>
		<true/>
		<key>com.apple.server.bluetooth.classic.xpc</key>
		<true/>
		<key>com.apple.server.bluetooth.general.xpc</key>
		<true/>
		<key>com.apple.server.bluetooth.le.att.xpc</key>
		<true/>
		<key>com.apple.server.bluetooth.le.pipe.xpc</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.ExposureNotification.UserNotification</key>
		<true/>
		<key>com.apple.wirelessproxd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/bluetoothd</string>
	</array>
	<key>PublishesEvents</key>
	<array>
		<string>com.apple.bluetooth.connections</string>
		<string>com.apple.bluetooth.discovery</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
