# beta/macOS/System/Library/LaunchDaemons/com.apple.AssetCacheLocatorService.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.AssetCacheLocatorService</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.AssetCacheLocatorService</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AssetCacheServices.framework/Versions/A/XPCServices/AssetCacheLocatorService.xpc/Contents/MacOS/AssetCacheLocatorService</string>
		<string>-d</string>
	</array>
	<key>UserName</key>
	<string>_assetcache</string>
</dict>
</plist>
