# beta/macOS/System/Library/LaunchDaemons/com.apple.akd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.akd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.ak.anisette.xpc</key>
		<true/>
		<key>com.apple.ak.auth.xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AuthKit.framework/Versions/A/Support/akd</string>
	</array>
</dict>
</plist>
