# beta/macOS/System/Library/LaunchDaemons/com.apple.nsurlsessiond.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>GroupName</key>
	<string>_nsurlsessiond</string>
	<key>Label</key>
	<string>com.apple.nsurlsessiond_privileged</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.nsurlsessiond</key>
		<true/>
		<key>com.apple.nsurlsessiond-launchondemand</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/nsurlsessiond</string>
		<string>--privileged</string>
	</array>
	<key>QueueDirectories</key>
	<array>
		<string>/var/root/Library/com.apple.nsurlsessiond</string>
	</array>
	<key>UserName</key>
	<string>_nsurlsessiond</string>
</dict>
</plist>
