# beta/macOS/System/Library/LaunchDaemons/com.apple.liquiddetectiond.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.apple.liquiddetectiond</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.liquiddetection.io-port-event</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IOPort</string>
			</dict>
			<key>com.apple.liquiddetection.ldcm-event</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyExistsMatch</key>
				<string>LiquidDetected</string>
				<key>IOProviderClass</key>
				<string>IOPortFeatureLDCM</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.liquiddetectionactivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>300</integer>
				<key>Interval</key>
				<integer>3600</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.ldcm.ldcm-server</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Background</string>
	<key>Program</key>
	<string>/System/Library/CoreServices/liquiddetectiond.app/liquiddetectiond</string>
	<key>ProgramArguments</key>
	<array/>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
