# beta/macOS/System/Library/LaunchDaemons/com.apple.cfnetwork.cfnetworkagent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>JavaScriptCoreUseJIT</key>
		<string>0</string>
	</dict>
	<key>Label</key>
	<string>com.apple.cfnetwork.cfnetworkagent</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.cfnetwork.cfnetworkagent</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/CFNetwork.framework/Versions/A/Support/CFNetworkAgent</string>
	</array>
</dict>
</plist>
