# beta/macOS/System/Library/LaunchDaemons/com.apple.memory-maintenance.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AbandonProcessGroup</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.memory-maintenance</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.memory-maintenance.ane-abandonment-check</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>3600</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.memory-maintenance.report-ecc-status</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>kBGSTProtectionClassD</string>
			</dict>
			<key>com.apple.memory-maintenance.report-jetsam-telemetry</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>28800</integer>
				</dict>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.trial.NamespaceUpdate.COREOS_FAST_PREWARMING</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.COREOS_FAST_PREWARMING</string>
			</dict>
			<key>com.apple.trial.NamespaceUpdate.COREOS_GMPOWER_VM_TUNING_PAGE_SHORTAGE_THRESHOLDS</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.COREOS_GMPOWER_VM_TUNING_PAGE_SHORTAGE_THRESHOLDS</string>
			</dict>
			<key>com.apple.trial.NamespaceUpdate.COREOS_VM_OBJECT_COPY_DELAYED_NO_WAIT</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.COREOS_VM_OBJECT_COPY_DELAYED_NO_WAIT</string>
			</dict>
			<key>com.apple.trial.NamespaceUpdate.FREEZER_POLICIES</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.FREEZER_POLICIES</string>
			</dict>
		</dict>
	</dict>
	<key>LowPriorityIO</key>
	<true/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mmaintenanced</key>
		<true/>
		<key>com.apple.mmaintenanced.server</key>
		<dict>
			<key>HostSpecialPort</key>
			<integer>33</integer>
		</dict>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/mmaintenanced</string>
	</array>
</dict>
</plist>
