# beta/macOS/System/Library/LaunchDaemons/com.apple.icloud.findmydeviced.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.icloud.findmydeviced</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.UserWasRemovedNotification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.UserWasRemovedNotification</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.ManagedClient.ActivationLockAllowedStateDidChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ManagedClient.ActivationLockAllowedStateDidChange</string>
			</dict>
			<key>com.apple.icloud.searchparty.selfbeaconchanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.searchparty.selfbeaconchanged</string>
			</dict>
			<key>com.apple.locationd/Prefs</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd/Prefs</string>
			</dict>
			<key>com.apple.system.hostname</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.hostname</string>
			</dict>
			<key>com.apple.system.loginwindow.desktopUp</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.loginwindow.desktopUp</string>
			</dict>
		</dict>
		<key>com.apple.nsurlsessiond-events</key>
		<dict>
			<key>LaunchOnDemand</key>
			<dict>
				<key>ClientIdentifier</key>
				<string>com.apple.icloud.findmydeviced</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.icloud.findmydeviced.launch-on-reboot</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>120</integer>
				<key>GracePeriod</key>
				<integer>30</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.icloud.findmydeviced</key>
		<true/>
		<key>com.apple.icloud.findmydeviced.aps-demo</key>
		<true/>
		<key>com.apple.icloud.findmydeviced.aps-development</key>
		<true/>
		<key>com.apple.icloud.findmydeviced.aps-production</key>
		<true/>
		<key>com.apple.icloud.findmydeviced.ua-services</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/findmydeviced</string>
	<key>Umask</key>
	<integer>23</integer>
</dict>
</plist>
