# beta/macOS/System/Library/LaunchDaemons/com.apple.dprivacyd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.dprivacyd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.DifferentialPrivacy.ReportFilesMaintenance</key>
			<dict>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.DifferentialPrivacy.StorageCulling</key>
			<dict>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.DifferentialPrivacy.StorageMaintenance</key>
			<dict>
				<key>Interval</key>
				<integer>43200</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.DifferentialPrivacy.TransparencyLogMaintenance</key>
			<dict>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/dprivacyd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/dprivacyd</string>
	</array>
</dict>
</plist>
