# beta/macOS/System/Library/LaunchDaemons/com.apple.accessoryd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.accessoryd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.accessoryd.matching.A2868.billboard</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOParentMatch</key>
				<dict>
					<key>IOPropertyMatch</key>
					<dict>
						<key>idProduct</key>
						<integer>5904</integer>
						<key>idVendor</key>
						<integer>1452</integer>
					</dict>
					<key>IOProviderClass</key>
					<string>IOUSBHostInterface</string>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleUSBHostBillboardDevice</string>
			</dict>
			<key>com.apple.accessoryd.matching.A2869.billboard</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOParentMatch</key>
				<dict>
					<key>IOPropertyMatch</key>
					<dict>
						<key>idProduct</key>
						<integer>5905</integer>
						<key>idVendor</key>
						<integer>1452</integer>
					</dict>
					<key>IOProviderClass</key>
					<string>IOUSBHostInterface</string>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleUSBHostBillboardDevice</string>
			</dict>
			<key>com.apple.accessoryd.matching.authCP</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>AppleAuthCP</string>
			</dict>
			<key>com.apple.accessoryd.matching.c99</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>Active</key>
					<true/>
					<key>AuthenticationRequired</key>
					<true/>
					<key>ParentPortType</key>
					<integer>17</integer>
					<key>TransportType</key>
					<integer>1</integer>
				</dict>
				<key>IOProviderClass</key>
				<string>IOPortTransportStateCC</string>
			</dict>
			<key>com.apple.accessoryd.matching.usbHost</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>SupportsIAP</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>IOUSBHostInterface</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.accessories.connection-info-server</key>
		<true/>
		<key>com.apple.accessories.externalaccessory-server</key>
		<true/>
		<key>com.apple.accessories.hid</key>
		<true/>
		<key>com.apple.accessories.now-playing</key>
		<true/>
		<key>com.apple.accessories.transport-server</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CoreAccessories.framework/Support/accessoryd</string>
	</array>
</dict>
</plist>
