# beta/macOS/System/Library/LaunchDaemons/com.apple.DesktopServicesHelper.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.DesktopServicesHelper</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.DesktopServicesHelper</key>
		<dict/>
	</dict>
	<key>MultipleInstances</key>
	<true/>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/DesktopServicesPriv.framework/Resources/DesktopServicesHelper</string>
</dict>
</plist>
