# beta/macOS/System/Library/LaunchDaemons/com.apple.remotemanagementd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>GroupName</key>
	<string>_rmd</string>
	<key>Label</key>
	<string>com.apple.remotemanagementd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.remotemanagement.on-reboot</key>
			<dict>
				<key>NonRepeatingTask</key>
				<dict>
					<key>ScheduleAfter</key>
					<integer>120</integer>
					<key>TrySchedulingBefore</key>
					<integer>150</integer>
				</dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequiresExternalPower</key>
				<false/>
			</dict>
			<key>com.apple.remotemanagement.periodic-sync</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>68400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.remotemanagement.cloudConfigAvailable</key>
			<dict>
				<key>Name</key>
				<string>CPCloudConfigIsMDMv2Notification</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>Accounts Changed Notification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.remotemanagement.accountsChanged</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.remotemanagementd.http.apns-dev</key>
		<true/>
		<key>com.apple.aps.remotemanagementd.http.apns-prod</key>
		<true/>
		<key>com.apple.remotemanagementd</key>
		<true/>
		<key>com.apple.remotemanagementd.store</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/RemoteManagement.framework/remotemanagementd</string>
	<key>PublishesEvents</key>
	<string>com.apple.remotemanagement.configs.daemon</string>
	<key>UserName</key>
	<string>_rmd</string>
</dict>
</plist>
