# beta/macOS/System/Library/LaunchDaemons/com.apple.PerfPowerServicesExtended.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>OS_ACTIVITY_MODE</key>
		<string>enable</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/var/db/powerlog/powerlogd.launchd</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.PerfPowerServicesExtended</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.powerlogd.XPCService.xpc</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/PerfPowerServicesExtended</string>
	</array>
</dict>
</plist>
