# beta/macOS/System/Library/LaunchDaemons/com.apple.AppleCredentialManagerDaemon.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.AppleCredentialManagerDaemon</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.fsevents.matching</key>
		<dict>
			<key>AppleSetupCookieChanged</key>
			<dict>
				<key>Path</key>
				<string>/var/db/.AppleSetupDone</string>
			</dict>
		</dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>ACMTRMEvent_CheckAppleSetup</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyExistsMatch</key>
				<string>ACMTRMEvent_CheckAppleSetup</string>
				<key>IOProviderClass</key>
				<string>AppleCredentialManager</string>
			</dict>
			<key>ACMTRMEvent_SaveData</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyExistsMatch</key>
				<string>ACMTRMEvent_SaveData</string>
				<key>IOProviderClass</key>
				<string>AppleCredentialManager</string>
			</dict>
			<key>ACMTRMEvent_ScheduleWakeup</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyExistsMatch</key>
				<string>ACMTRMEvent_ScheduleWakeup</string>
				<key>IOProviderClass</key>
				<string>AppleCredentialManager</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.MCX._managementStatusChangedForDomains</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MCX._managementStatusChangedForDomains</string>
			</dict>
			<key>com.apple.private.restrict-post.fmip.lostmode.enable</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.private.restrict-post.fmip.lostmode.enable</string>
			</dict>
		</dict>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AppleCredentialManager.framework/AppleCredentialManagerDaemon</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
