# beta/macOS/System/Library/LaunchDaemons/com.apple.PowerUIAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.PowerUIAgent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.powerlog.batteryServiceNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.powerlog.batteryServiceNotification</string>
			</dict>
			<key>com.apple.powerui.requiredFullCharge</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.powerui.requiredFullCharge</string>
			</dict>
			<key>com.apple.powerui.smartcharge</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.powerui.smartcharge</string>
			</dict>
			<key>com.apple.smartcharging.defaultschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.smartcharging.defaultschanged</string>
			</dict>
			<key>com.apple.system.powermanagement.useractivity2</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powermanagement.useractivity2</string>
			</dict>
			<key>com.apple.system.powersources.source</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powersources.source</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.poweruiagent.reportAnalyticsRepeating</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>28800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.poweruiagent.reportMonthlyAnalytics</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>604800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
			<key>com.apple.poweruiagent.runOnBoot</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>0</integer>
				<key>GracePeriod</key>
				<integer>0</integer>
				<key>Priority</key>
				<string>Utility</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.powerui.bdcdata</key>
		<true/>
		<key>com.apple.powerui.smartChargeManager</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.powerui.lowpowermode</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.powerui.smartcharging</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/usr/libexec/PowerUIAgent</string>
</dict>
</plist>
