# beta/macOS/System/Library/LaunchDaemons/com.apple.audio.systemsoundserverd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.audio.systemsoundserverd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.audio.SystemSoundServer-OSX</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>CoreAudioServices/sss_in_audiomxd</string>
			<key>#Then</key>
			<true/>
		</dict>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/systemsoundserverd</string>
	</array>
</dict>
</plist>
