# beta/macOS/System/Library/LaunchDaemons/com.apple.mediaremoted.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.apple.mediaremoted</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.AirTunes.DACP.device-prevent-playback</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.device-prevent-playback</string>
			</dict>
			<key>com.apple.AirTunes.DACP.devicevolume</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.devicevolume</string>
			</dict>
			<key>com.apple.AirTunes.DACP.devicevolumechanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.devicevolumechanged</string>
			</dict>
			<key>com.apple.AirTunes.DACP.mutetoggle</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.mutetoggle</string>
			</dict>
			<key>com.apple.AirTunes.DACP.nextitem</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.nextitem</string>
			</dict>
			<key>com.apple.AirTunes.DACP.pause</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.pause</string>
			</dict>
			<key>com.apple.AirTunes.DACP.play</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.play</string>
			</dict>
			<key>com.apple.AirTunes.DACP.previtem</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.previtem</string>
			</dict>
			<key>com.apple.AirTunes.DACP.repeatadv</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.repeatadv</string>
			</dict>
			<key>com.apple.AirTunes.DACP.shuffletoggle</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.shuffletoggle</string>
			</dict>
			<key>com.apple.AirTunes.DACP.volumedown</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.volumedown</string>
			</dict>
			<key>com.apple.AirTunes.DACP.volumeup</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AirTunes.DACP.volumeup</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mediacontrol.xpc</key>
		<true/>
		<key>com.apple.mediaremoted.xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/MediaRemote.framework/Support/mediaremoted</string>
	</array>
</dict>
</plist>
