# diff/iOS/System/Library/NanoLaunchDaemons/com.apple.nanonewscd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.nanonewscd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.nanonews.service.companion</key>
		<true/>
		<key>com.apple.private.alloy.news-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/NewsServices.framework/nanonewscd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
