# diff/iOS/System/Library/NanoLaunchDaemons/com.apple.companionmessagesd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>Label</key>
	<string>com.apple.companionmessagesd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.TextInput.RemoveAllDynamicDictionariesNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TextInput.RemoveAllDynamicDictionariesNotification</string>
			</dict>
			<key>com.apple.messages.idslaunchnotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.messages.idslaunchnotification</string>
			</dict>
			<key>com.apple.quickboard.classa.idslaunchnotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.quickboard.classa.idslaunchnotification</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.companionmessagesd.DeviceSaltSync</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresClassA</key>
				<true/>
			</dict>
			<key>com.apple.companionmessagesd.ResponseKitSync</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.pairedsync.quickboard</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/ChatKit.framework/companionmessagesd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
