# diff/iOS/System/Library/NanoLaunchDaemons/com.apple.companioninferenceserviced.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.companioninferenceserviced</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.companion-inference</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.companion-inference</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.companioninferenceserviced</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CompanionInferenceCore.framework/companioninferenceserviced</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
