# beta/Library/LaunchDaemons/com.apple.anomalydetectiond.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.anomalydetectiond</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.anomalydetectiond.CSFolderMonitor.encrypted</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>60</integer>
				<key>Interval</key>
				<real>14400</real>
				<key>NetworkTransferDirection</key>
				<string>Upload</string>
				<key>NetworkUploadSize</key>
				<integer>500000</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.anomalydetectiond.CSFolderMonitor.metadata</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>CPUIntensive</key>
				<true/>
				<key>Delay</key>
				<integer>60</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>NetworkDownloadSize</key>
				<integer>1000</integer>
				<key>NetworkTransferDirection</key>
				<string>Download</string>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequiresClassB</key>
				<true/>
			</dict>
			<key>com.apple.anomalydetectiond.CSFolderMonitor.protodata</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>60</integer>
				<key>Interval</key>
				<real>14400</real>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresClassB</key>
				<true/>
			</dict>
			<key>com.apple.anomalydetectiond.CSFolderMonitorIgneous.metadata</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>60</integer>
				<key>Interval</key>
				<real>14400</real>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresClassB</key>
				<true/>
			</dict>
			<key>com.apple.anomalydetectiond.CSFolderMonitorIgneous.upload</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>CPUIntensive</key>
				<true/>
				<key>Delay</key>
				<integer>60</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>NetworkDownloadSize</key>
				<integer>1000</integer>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>NetworkUploadSize</key>
				<string>500000</string>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequiresClassB</key>
				<true/>
			</dict>
			<key>com.apple.anomalydetectiond.CoreAnalytics.UserInfo</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToHardware</key>
	<dict>
		<key>targettype</key>
		<array>
			<string>D27</string>
			<string>D28</string>
			<string>D37</string>
			<string>D38</string>
			<string>D47</string>
			<string>D48</string>
			<string>D73</string>
			<string>D74</string>
			<string>D83</string>
			<string>D84</string>
			<string>D93</string>
			<string>D94</string>
			<string>V59</string>
			<string>D23</string>
			<string>V53</string>
			<string>V54</string>
			<string>V57</string>
			<string>N104</string>
			<string>D421</string>
			<string>D431</string>
			<string>D49</string>
			<string>D52g</string>
			<string>D53p</string>
			<string>D53g</string>
			<string>D54p</string>
			<string>D16</string>
			<string>D17</string>
			<string>D63</string>
			<string>D64</string>
		</array>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.anomalydetectiond.kappa</key>
		<true/>
		<key>com.apple.anomalydetectiond.trigger</key>
		<true/>
		<key>com.apple.private.alloy.kappa-idswake</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/anomalydetectiond</string>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
	<key>_Conclave</key>
	<string>com.apple.anomalydetectiond.conclave</string>
</dict>
</plist>
