# beta/Library/LaunchDaemons/com.apple.ReportMemoryException.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>__CFPREFERENCES_AVOID_DAEMON</key>
		<string>1</string>
	</dict>
	<key>GroupName</key>
	<string>_reportmemoryexception</string>
	<key>Label</key>
	<string>com.apple.ReportMemoryException</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.ReportMemoryException</key>
		<dict>
			<key>DrainMessagesOnCrash</key>
			<string>All</string>
		</dict>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/ReportMemoryException</string>
	</array>
	<key>UserName</key>
	<string>_reportmemoryexception</string>
</dict>
</plist>
