# beta/Library/LaunchDaemons/com.apple.CommCenterMobileHelper.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>CFBundleIdentifier</key>
	<string>com.apple.CommCenterMobileHelper</string>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>GroupName</key>
	<string>mobile</string>
	<key>Label</key>
	<string>com.apple.CommCenterMobileHelper</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.ManagedConfiguration.profileListChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ManagedConfiguration.profileListChanged</string>
			</dict>
			<key>com.apple.purplebuddy.setupdone</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.purplebuddy.setupdone</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.abm.helper.mobile</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.commcenter.mobile-helper</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.commcenter.mobile-helper-rest</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.commcenter.mobile-helper-xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/CoreTelephony.framework/Support/CommCenterMobileHelper</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>0</integer>
	<key>UserName</key>
	<string>mobile</string>
	<key>application-identifier</key>
	<string>com.apple.CommCenterMobileHelper</string>
	<key>com.apple.private.attribution.implicitly-assumed-identity</key>
	<dict>
		<key>type</key>
		<string>path</string>
		<key>value</key>
		<string>/System/Library/Frameworks/CoreTelephony.framework/Support/CommCenterMobileHelper</string>
	</dict>
</dict>
</plist>
