# beta/Library/LaunchDaemons/com.apple.devicemanagementclient.teslad.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.devicemanagementclient.teslad</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.managedconfiguration.teslad.deppush-prod</key>
		<true/>
		<key>com.apple.managedconfiguration.teslad</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/teslad</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
