# beta/Library/LaunchDaemons/com.apple.Maps.mapssyncd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.Maps.mapssyncd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.Maps.MapsSync.activity</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.Maps.MapsSync.migration</key>
			<dict>
				<key>NonRepeatingTask</key>
				<dict/>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RequiresBuddyComplete</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.Maps.MapsSync.cache</key>
		<true/>
		<key>com.apple.Maps.MapsSync.service</key>
		<true/>
		<key>com.apple.Maps.MapsSync.store</key>
		<true/>
		<key>com.apple.Maps.MapsSync.xpc</key>
		<true/>
		<key>com.apple.aps.mapssyncd</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/MapsSync.framework/mapssyncd</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
