# beta/Library/LaunchDaemons/com.apple.UsageTrackingAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.UsageTrackingAgent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.ManagedSettings.effective-settings.changed</key>
		<dict>
			<key>com.apple.UsageTrackingAgent.effective-settings.changed</key>
			<dict>
				<key>SettingsGroups</key>
				<array>
					<string>deviceActivity</string>
				</array>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.UsageTrackingAgent.distributed-notification.duet-application-usage-did-tombstone</key>
			<dict>
				<key>Name</key>
				<string>_DKKnowledgeStorageDidTombstoneEventsNotification</string>
				<key>Object</key>
				<string>/app/usage</string>
			</dict>
			<key>com.apple.UsageTrackingAgent.distributed-notification.duet-now-playing-usage-did-tombstone</key>
			<dict>
				<key>Name</key>
				<string>_DKKnowledgeStorageDidTombstoneEventsNotification</string>
				<key>Object</key>
				<string>/media/nowPlaying</string>
			</dict>
			<key>com.apple.UsageTrackingAgent.distributed-notification.duet-sync-coordinator-did-idle</key>
			<dict>
				<key>Name</key>
				<string>_DKSync2CoordinatorDidIdleNotification</string>
			</dict>
			<key>com.apple.UsageTrackingAgent.distributed-notification.duet-video-usage-did-tombstone</key>
			<dict>
				<key>Name</key>
				<string>_DKKnowledgeStorageDidTombstoneEventsNotification</string>
				<key>Object</key>
				<string>/app/mediaUsage</string>
			</dict>
			<key>com.apple.UsageTrackingAgent.distributed-notification.duet-web-domain-usage-did-tombstone</key>
			<dict>
				<key>Name</key>
				<string>_DKKnowledgeStorageDidTombstoneEventsNotification</string>
				<key>Object</key>
				<string>/app/webUsage</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>SignificantTimeChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>SignificantTimeChangeNotification</string>
			</dict>
			<key>com.apple.UsageTrackingAgent.registration.application</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.UsageTrackingAgent.registration.application</string>
			</dict>
			<key>com.apple.UsageTrackingAgent.registration.now-playing</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.UsageTrackingAgent.registration.now-playing</string>
			</dict>
			<key>com.apple.UsageTrackingAgent.registration.video</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.UsageTrackingAgent.registration.video</string>
			</dict>
			<key>com.apple.UsageTrackingAgent.registration.web-domain</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.UsageTrackingAgent.registration.web-domain</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.UsageTrackingAgent.activity.device-activity-collection</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>600</integer>
				<key>Interval</key>
				<integer>3600</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.UsageTrackingAgent.activity.reset-device-activity-alarms</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>60</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequiresClassA</key>
				<true/>
			</dict>
			<key>com.apple.UsageTrackingAgent.activity.uncategorized-web-usage</key>
			<dict>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.UsageTrackingAgent</key>
		<true/>
		<key>com.apple.UsageTrackingAgent.private</key>
		<true/>
		<key>com.apple.aps.UsageTrackingAgent</key>
		<true/>
		<key>com.apple.private.alloy.usagetracking-idswake</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/UsageTracking.framework/UsageTrackingAgent</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
