# beta/Library/LaunchDaemons/com.apple.sleepd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>OS_ACTIVITY_ENABLE_DYNAMIC</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.sleepd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.alarm</key>
		<true/>
		<key>com.apple.coreduetcontext.client_event_stream</key>
		<dict>
			<key>notification</key>
			<dict>
				<key>client-identifier</key>
				<string>com.apple.sleepd.contextstore-registration</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>Application Installed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>Application Uninstalled</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CSLDisableWristDetectionChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>CSLDisableWristDetectionChangedNotification</string>
			</dict>
			<key>SignificantTimeChange</key>
			<dict>
				<key>Notification</key>
				<string>SignificantTimeChangeNotification</string>
			</dict>
			<key>com.apple.Carousel.wristStateChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.Carousel.wristStateChanged</string>
			</dict>
			<key>com.apple.healthlite.SleepDetectedActivity</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.healthlite.SleepDetectedActivity</string>
			</dict>
			<key>com.apple.healthlite.SleepSessionEndRequest</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.healthlite.SleepSessionEndRequest</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.sleepd</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.sleepd</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.nanoregistry.watchdidbecomeactive</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.watchdidbecomeactive</string>
			</dict>
			<key>com.apple.sleep.sync.SleepRecordDidChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sleep.sync.SleepRecordDidChange</string>
			</dict>
			<key>com.apple.sleep.sync.SleepScheduleDidChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sleep.sync.SleepScheduleDidChange</string>
			</dict>
			<key>com.apple.sleep.sync.SleepSettingsDidChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sleep.sync.SleepSettingsDidChange</string>
			</dict>
			<key>com.apple.sleepd.analytics</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sleepd.analytics</string>
			</dict>
			<key>com.apple.sleepd.cloudkit.reset</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sleepd.cloudkit.reset</string>
			</dict>
			<key>com.apple.sleepd.diagnostics</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sleepd.diagnostics</string>
			</dict>
			<key>com.apple.sleepd.ids.test</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sleepd.ids.test</string>
			</dict>
			<key>com.apple.system.powersources.source</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powersources.source</string>
			</dict>
			<key>com.apple.system.timezone</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.timezone</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.private.alloy.sleep.classd-idswake</key>
		<true/>
		<key>com.apple.sleepd.sleepserver</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/SleepDaemon.framework/sleepd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
