# beta/Library/LaunchDaemons/com.apple.biometrickitd.pearl.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LLVM_PROFILE_FILE</key>
		<string>/tmp/biometrickitd.%p.profraw</string>
	</dict>
	<key>Label</key>
	<string>com.apple.biometrickitd.pearl</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.pearld</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>ApplePearlSEPDriver</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.managedconfiguration.passcodechanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.passcodechanged</string>
			</dict>
			<key>com.apple.mobile.keybagd.lock_status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.lock_status</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.pearld</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/biometrickitd</string>
	</array>
	<key>_LimitLoadToDeviceTree</key>
	<string>pearl-sep</string>
</dict>
</plist>
