# beta/Library/LaunchDaemons/com.apple.deleted_helper.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.deleted_helper</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.dispatch.vfs</key>
		<dict>
			<key>Monitor Low Disk Conditions</key>
			<dict>
				<key>DesiredDisk</key>
				<true/>
				<key>LowDisk</key>
				<integer>0</integer>
				<key>NearLowDisk</key>
				<true/>
				<key>VeryLowDisk</key>
				<integer>0</integer>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.deleted_helper.daily</key>
			<dict>
				<key>Interval</key>
				<integer>86400</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.cache_delete_orphan_dir_handler</key>
		<true/>
		<key>com.apple.deleted_helper</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CacheDelete.framework/deleted_helper</string>
	</array>
</dict>
</plist>
