# beta/Library/LaunchDaemons/com.apple.accountsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamMemoryLimit</key>
		<integer>2000</integer>
	</dict>
	<key>Label</key>
	<string>com.apple.accountsd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.accountsd.accountmanager</key>
		<true/>
		<key>com.apple.accountsd.oauthsigner</key>
		<true/>
		<key>com.apple.accountsd.oopa</key>
		<true/>
		<key>com.apple.private.alloy.accountssync-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/Accounts.framework/accountsd</string>
	</array>
	<key>PublishesEvents</key>
	<dict>
		<key>com.apple.accounts.changes</key>
		<dict>
			<key>RequireTrusted</key>
			<true/>
		</dict>
	</dict>
	<key>ThrottleInterval</key>
	<integer>3</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
