# beta/Library/LaunchDaemons/com.apple.fitcore.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MallocStackLoggingDirectory</key>
		<string>/var/mobile/Library/Seymour</string>
		<key>NSRunningFromLaunchd</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.fitcore</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.authkit.trusted-device-list-changed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.authkit.trusted-device-list-changed</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.ManagedConfiguration.profileListChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ManagedConfiguration.profileListChanged</string>
			</dict>
			<key>com.apple.appstored.ActivitySubEntitlementsCacheUpdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appstored.ActivitySubEntitlementsCacheUpdated</string>
			</dict>
			<key>com.apple.fitness.FitnessAppInstalled</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.fitness.FitnessAppInstalled</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>Seymour</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.Seymour</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.Seymour.Scheduling.PostInstallActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.Seymour.Scheduling.PrivateDataActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<real>14400</real>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresClassB</key>
				<true/>
			</dict>
			<key>com.apple.Seymour.Scheduling.PrivateDataNetworkActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<real>14400</real>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassB</key>
				<true/>
			</dict>
			<key>com.apple.Seymour.Scheduling.PublicDataNetworkActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>1800</integer>
				<key>Interval</key>
				<integer>14400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.Seymour.Scheduling.PublicDataNetworkPowerActivity</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>GracePeriod</key>
				<real>14400</real>
				<key>Interval</key>
				<integer>28800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.Seymour.Scheduling.PublicDataPowerActivity</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>GracePeriod</key>
				<real>14400</real>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.Seymour.Scheduling.RandomPublicDataNetworkActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<real>14400</real>
				<key>Interval</key>
				<integer>28800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<integer>28800</integer>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.Seymour.Scheduling.RandomizedDailyPublicDataNetworkActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<integer>86400</integer>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.seymour</key>
		<true/>
		<key>com.apple.seymour</key>
		<true/>
		<key>com.apple.seymour.cache-delete</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/SeymourServices.framework/fitcored</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
