# beta/Library/LaunchDaemons/com.apple.FamilyControlsAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.FamilyControlsAgent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.family.family_updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.family.family_updated</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.FamilyControlsAgent</key>
		<true/>
		<key>com.apple.FamilyControlsAgent.private</key>
		<true/>
		<key>com.apple.aps.FamilyControlsAgent</key>
		<true/>
		<key>com.apple.private.alloy.familycontrols-idswake</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/Frameworks/FamilyControls.framework/FamilyControlsAgent</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
