# beta/Library/LaunchDaemons/com.apple.voicemail.vmd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>NSRunningFromLaunchd</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.voicemail.vmd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.CTTelephonyCenter</key>
		<dict>
			<key>kCTIndicatorsVoiceMailNotification</key>
			<dict>
				<key>Name</key>
				<string>kCTIndicatorsVoiceMailNotification</string>
			</dict>
			<key>kCTSettingRequestSuccessNotification</key>
			<dict>
				<key>Name</key>
				<string>kCTSettingRequestSuccessNotification</string>
			</dict>
			<key>kCTVoicemailDatabaseResetNotification</key>
			<dict>
				<key>Name</key>
				<string>kCTVoicemailDatabaseResetNotification</string>
			</dict>
			<key>kCTVoicemailWakeupNotification</key>
			<dict>
				<key>Name</key>
				<string>kCTVoicemailWakeupNotification</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>NewCarrierNotification</key>
			<dict>
				<key>Notification</key>
				<string>NewCarrierNotification</string>
			</dict>
			<key>NewOperatorNotification</key>
			<dict>
				<key>Notification</key>
				<string>NewOperatorNotification</string>
			</dict>
			<key>VMStoreSetTokenNotification</key>
			<dict>
				<key>Notification</key>
				<string>VMStoreSetTokenNotification</string>
			</dict>
			<key>VVMessageWaitingFallbackNotification</key>
			<dict>
				<key>Notification</key>
				<string>VVMessageWaitingFallbackNotification</string>
			</dict>
			<key>com.apple.commcenter.InternationalRoamingEDGE.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.commcenter.InternationalRoamingEDGE.changed</string>
			</dict>
			<key>com.apple.telephonyutilities.callservicesdaemon.voicemailcallended</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.telephonyutilities.callservicesdaemon.voicemailcallended</string>
			</dict>
			<key>com.apple.voicemail.ReloadService</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicemail.ReloadService</string>
			</dict>
			<key>com.apple.voicemail.VVVerifierCheckpointDictionaryChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicemail.VVVerifierCheckpointDictionaryChanged</string>
			</dict>
			<key>com.apple.voicemail.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicemail.changed</string>
			</dict>
			<key>kVMVoicemailTranscriptionTaskTranscribeAllVoicemails</key>
			<dict>
				<key>Notification</key>
				<string>kVMVoicemailTranscriptionTaskTranscribeAllVoicemails</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.voicemail.vmd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/VisualVoicemail.framework/vmd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>WorkingDirectory</key>
	<string>/tmp</string>
</dict>
</plist>
