# beta/Library/LaunchDaemons/com.apple.mobileactivationd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.mobileactivationd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.os-eligibility-domain.input-needed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.os-eligibility-domain.input-needed</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.mobileactivationd.DCRT.OOB</key>
			<dict>
				<key>Interval</key>
				<real>86400</real>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.mobileactivationd.DCRT.OOB.LoadSpreading</key>
			<dict>
				<key>Delay</key>
				<integer>0</integer>
				<key>Interval</key>
				<integer>2246400</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>2246400</integer>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<false/>
			</dict>
			<key>com.apple.mobileactivationd.SDCRT.OOB</key>
			<dict>
				<key>Interval</key>
				<real>86400</real>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.mobileactivationd.SDCRT.OOB.LoadSpreading</key>
			<dict>
				<key>Delay</key>
				<integer>0</integer>
				<key>Interval</key>
				<integer>2246400</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>2246400</integer>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<false/>
			</dict>
			<key>com.apple.mobileactivationd.splunklogging</key>
			<dict>
				<key>Interval</key>
				<real>86400</real>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.mobileactivationd.splunklogging.daily</key>
			<dict>
				<key>Interval</key>
				<real>86400</real>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.mobileactivationd.splunklogging.rollidentifier</key>
			<dict>
				<key>Interval</key>
				<integer>1209600</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mobileactivationd</key>
		<true/>
		<key>com.apple.mobileactivationd.lockdown</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/mobileactivationd</string>
	</array>
	<key>UserName</key>
	<string>root</string>
</dict>
</plist>
