# beta/Library/LaunchDaemons/com.apple.CarPlayApp.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>CFBundleIdentifier</key>
	<string>com.apple.CarPlayApp</string>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.CarPlayApp</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.CarPlayApp.cluster-theme-service</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.CarPlayApp.non-launching-service</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.CarPlayApp.oswatchdog</key>
		<true/>
		<key>com.apple.CarPlayApp.punch-through-service</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.CarPlayApp.service</key>
		<true/>
		<key>com.apple.CarPlayApp.status-bar-service</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.CarPlayApp.user-alerts-service</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.CarPlayApp.volume-notification-service</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.carkit.layer-metadata</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/CarPlay.app/CarPlay</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>_AdditionalProperties</key>
	<dict>
		<key>RunningBoard</key>
		<dict>
			<key>Managed</key>
			<false/>
			<key>Reported</key>
			<true/>
		</dict>
	</dict>
</dict>
</plist>
