# beta/visionOS/System/Library/LaunchDaemons/com.apple.StatusKitAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.StatusKitAgent</string>
	<key>LaunchEvents</key>
	<dict/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.StatusKit.local</key>
		<true/>
		<key>com.apple.StatusKit.presence</key>
		<true/>
		<key>com.apple.StatusKit.publish</key>
		<true/>
		<key>com.apple.StatusKit.subscribe</key>
		<true/>
		<key>com.apple.aps.StatusKit</key>
		<true/>
		<key>com.apple.aps.StatusKit.CloudKitMirroring</key>
		<true/>
		<key>com.apple.private.alloy.status.keysharing-idswake</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/StatusKit.framework/StatusKitAgent</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
