# beta/visionOS/System/Library/LaunchDaemons/com.apple.arkitd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.apple.arkitd</string>
	<key>LimitLoadToHardware</key>
	<dict>
		<key>chiprole</key>
		<array>
			<integer>0</integer>
		</array>
		<key>targettype</key>
		<array>
			<string>N301</string>
			<string>N301a</string>
		</array>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.CoreIKDebuggingMachService</key>
		<true/>
		<key>com.apple.arkit.service.accessoryTracking</key>
		<dict>
			<key>#Else</key>
			<true/>
			<key>#IfFeatureFlagEnabled</key>
			<string>ARKitCore/accessory_tracking_uses_peer_daemon</string>
			<key>#Then</key>
			<false/>
		</dict>
		<key>com.apple.arkit.service.acousticMesh</key>
		<true/>
		<key>com.apple.arkit.service.appClipCode</key>
		<true/>
		<key>com.apple.arkit.service.appleDeviceTracking</key>
		<true/>
		<key>com.apple.arkit.service.authorization</key>
		<true/>
		<key>com.apple.arkit.service.barcodeDetection</key>
		<true/>
		<key>com.apple.arkit.service.cameraFrameProvider</key>
		<true/>
		<key>com.apple.arkit.service.clipOn</key>
		<true/>
		<key>com.apple.arkit.service.collaborativeMapping</key>
		<true/>
		<key>com.apple.arkit.service.control</key>
		<true/>
		<key>com.apple.arkit.service.deflicker</key>
		<true/>
		<key>com.apple.arkit.service.depthInfieldCalibration</key>
		<true/>
		<key>com.apple.arkit.service.enhancedPassthrough</key>
		<true/>
		<key>com.apple.arkit.service.environmentLightEstimation</key>
		<true/>
		<key>com.apple.arkit.service.faceRepresentation</key>
		<true/>
		<key>com.apple.arkit.service.fastPathInProcessPrediction</key>
		<true/>
		<key>com.apple.arkit.service.focusRegion</key>
		<true/>
		<key>com.apple.arkit.service.frameDebug</key>
		<true/>
		<key>com.apple.arkit.service.gazeCoaching</key>
		<true/>
		<key>com.apple.arkit.service.gazeEnrollment</key>
		<true/>
		<key>com.apple.arkit.service.gazeTracking</key>
		<true/>
		<key>com.apple.arkit.service.handEnrollment</key>
		<true/>
		<key>com.apple.arkit.service.handTracking</key>
		<true/>
		<key>com.apple.arkit.service.imageTracking</key>
		<true/>
		<key>com.apple.arkit.service.lightMonitoring</key>
		<true/>
		<key>com.apple.arkit.service.materialSegmentation</key>
		<true/>
		<key>com.apple.arkit.service.mcamInfieldCalibration</key>
		<true/>
		<key>com.apple.arkit.service.meshReconstruction</key>
		<true/>
		<key>com.apple.arkit.service.metricDepth</key>
		<true/>
		<key>com.apple.arkit.service.motionCapture</key>
		<true/>
		<key>com.apple.arkit.service.networkingSynchronization</key>
		<true/>
		<key>com.apple.arkit.service.objectTracking</key>
		<dict>
			<key>#Else</key>
			<true/>
			<key>#IfFeatureFlagEnabled</key>
			<string>ARKitCore/object_tracking_uses_peer_daemon</string>
			<key>#Then</key>
			<false/>
		</dict>
		<key>com.apple.arkit.service.peerSynchronization</key>
		<true/>
		<key>com.apple.arkit.service.personPetDetection</key>
		<true/>
		<key>com.apple.arkit.service.planeDetection</key>
		<true/>
		<key>com.apple.arkit.service.prediction</key>
		<true/>
		<key>com.apple.arkit.service.prototyping</key>
		<true/>
		<key>com.apple.arkit.service.renderProperties</key>
		<true/>
		<key>com.apple.arkit.service.respirationTracking</key>
		<true/>
		<key>com.apple.arkit.service.roomTracking</key>
		<true/>
		<key>com.apple.arkit.service.sceneClassification</key>
		<true/>
		<key>com.apple.arkit.service.sceneGeometry</key>
		<true/>
		<key>com.apple.arkit.service.sceneObjectDetection2D</key>
		<true/>
		<key>com.apple.arkit.service.sceneObjectDetection3D</key>
		<true/>
		<key>com.apple.arkit.service.sceneUnderstanding</key>
		<true/>
		<key>com.apple.arkit.service.secondaryRepresentation</key>
		<true/>
		<key>com.apple.arkit.service.sharedCoordinateSpace</key>
		<true/>
		<key>com.apple.arkit.service.spatialSceneClassification</key>
		<true/>
		<key>com.apple.arkit.service.systemModeControl</key>
		<true/>
		<key>com.apple.arkit.service.userBodyEnrollment</key>
		<true/>
		<key>com.apple.arkit.service.userBodyTracking</key>
		<true/>
		<key>com.apple.arkit.service.userBodyTrackingLite</key>
		<true/>
		<key>com.apple.arkit.service.visualDepth</key>
		<true/>
		<key>com.apple.arkit.service.visualLogger</key>
		<true/>
		<key>com.apple.arkit.service.watchdog</key>
		<true/>
		<key>com.apple.arkit.service.worldTracking</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/arkitd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
	<key>_PanicOnCrash</key>
	<dict>
		<key>#Else</key>
		<false/>
		<key>#IfFeatureFlagEnabled</key>
		<string>ARKitCore/arkitd_panic_on_crash</string>
		<key>#Then</key>
		<dict>
			<key>PanicOnConsecutiveCrash</key>
			<true/>
			<key>PanicOnCrashDeadline</key>
			<integer>60</integer>
		</dict>
	</dict>
</dict>
</plist>
