# beta/visionOS/System/Library/LaunchDaemons/com.apple.midiserver-ios.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.apple.midiserver</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bluetooth.connections</key>
		<dict>
			<key>MyMIDIService</key>
			<dict>
				<key>serviceUUIDs</key>
				<array>
					<string>03B80E5A-EDE8-4B33-A751-6CE34EC4C700</string>
				</array>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.midiserver</key>
		<true/>
		<key>com.apple.midiserver.io</key>
		<dict>
			<key>HideUntilCheckIn</key>
			<true/>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/CoreMIDI.framework/MIDIServer</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
