# beta/visionOS/System/Library/LaunchDaemons/com.apple.softwareupdateservicesd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressureExit</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/var/mobile/Library/SoftwareUpdate/susdAUKeepAlive</key>
			<true/>
			<key>/var/mobile/Library/SoftwareUpdate/susdKeepAlive</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.softwareupdateservicesd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.MobileSoftwareUpdate.OSVersionChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileSoftwareUpdate.OSVersionChanged</string>
			</dict>
			<key>com.apple.softwareupdate.autoinstall.startInstall</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdate.autoinstall.startInstall</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.SUCoreConfigScheduledScan</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.SUCoreConfigScheduledScan</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.activity.autoDownload</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.activity.autoDownload</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.activity.autoDownloadEnd</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.activity.autoDownloadEnd</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.activity.autoInstallEnd</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.activity.autoInstallEnd</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.activity.autoInstallUnlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.activity.autoInstallUnlock</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.activity.autoScan</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.activity.autoScan</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.activity.delayEndScan</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.activity.delayEndScan</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.activity.emergencyAutoScan</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.activity.emergencyAutoScan</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.activity.installAlert</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.activity.installAlert</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.activity.presentBanner</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.activity.presentBanner</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.activity.rollbackReboot</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.activity.rollbackReboot</string>
			</dict>
			<key>com.apple.softwareupdateservicesd.activity.splatAutoScan</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.softwareupdateservicesd.activity.splatAutoScan</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.softwareupdateservicesd</key>
		<true/>
		<key>com.apple.softwareupdateservicesui.controller</key>
		<true/>
		<key>com.apple.softwareupdateservicesui.tool</key>
		<true/>
	</dict>
	<key>OnDemand</key>
	<true/>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/SoftwareUpdateServices.framework/Support/softwareupdateservicesd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>86400</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
