# beta/visionOS/System/Library/LaunchDaemons/com.apple.PerfPowerServices.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.PerfPowerServices</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.powerlogd.dailyTasks</key>
			<dict>
				<key>Delay</key>
				<integer>79200</integer>
				<key>GracePeriod</key>
				<integer>7200</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.powerlog-GRRRR-Lite</key>
		<true/>
		<key>com.apple.powerlog.plxpclogger.xpc</key>
		<true/>
		<key>com.apple.powerlog_mover</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/PerfPowerServices</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
