# beta/visionOS/System/Library/LaunchDaemons/com.apple.SurfBoard.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EmbeddedHomeScreen</key>
	<true/>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.SurfBoard</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>RemoteCompositorServer</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.RemoteCompositorServer.RenderSession</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadFromHardware</key>
	<dict>
		<key>osenvironment</key>
		<array>
			<string>diagnostics</string>
		</array>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>PurpleSystemAppPort</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.DragUI.druid.system</key>
		<true/>
		<key>com.apple.UIKit.KeyboardManagement.hosted</key>
		<true/>
		<key>com.apple.UIKit.statusbarserver</key>
		<true/>
		<key>com.apple.accessibility.AXSurfBoardServer</key>
		<true/>
		<key>com.apple.aps.usernotifications</key>
		<true/>
		<key>com.apple.bulletinboard.dataproviderconnection</key>
		<true/>
		<key>com.apple.bulletinboard.observerconnection</key>
		<true/>
		<key>com.apple.bulletinboard.serverconduitconnection</key>
		<true/>
		<key>com.apple.bulletinboard.settingsconnection</key>
		<true/>
		<key>com.apple.corespotlight.daemon.notifications</key>
		<true/>
		<key>com.apple.devicesharing.guestusermodeservice</key>
		<true/>
		<key>com.apple.frontboard.systemappservices</key>
		<true/>
		<key>com.apple.surfboard.appledevicetrackingservice</key>
		<true/>
		<key>com.apple.surfboard.applicationservice</key>
		<true/>
		<key>com.apple.surfboard.authenticationservice</key>
		<true/>
		<key>com.apple.surfboard.breakthroughservice</key>
		<true/>
		<key>com.apple.surfboard.compositorservice</key>
		<true/>
		<key>com.apple.surfboard.debugservice</key>
		<true/>
		<key>com.apple.surfboard.deviceresetservice</key>
		<true/>
		<key>com.apple.surfboard.displayadjustmentservice</key>
		<true/>
		<key>com.apple.surfboard.draganddropservice</key>
		<true/>
		<key>com.apple.surfboard.entityaccessservice</key>
		<true/>
		<key>com.apple.surfboard.entityinteractionservice</key>
		<true/>
		<key>com.apple.surfboard.eyesightservice</key>
		<true/>
		<key>com.apple.surfboard.hardwareassertionservice</key>
		<true/>
		<key>com.apple.surfboard.immersionservice</key>
		<true/>
		<key>com.apple.surfboard.keyboardservice</key>
		<true/>
		<key>com.apple.surfboard.launcherservice</key>
		<true/>
		<key>com.apple.surfboard.lockscreenservice</key>
		<true/>
		<key>com.apple.surfboard.oswatchdog</key>
		<true/>
		<key>com.apple.surfboard.placementservice</key>
		<true/>
		<key>com.apple.surfboard.scene-resizing-service</key>
		<true/>
		<key>com.apple.surfboard.scenesessionservice</key>
		<true/>
		<key>com.apple.surfboard.scenesessionuiservice</key>
		<true/>
		<key>com.apple.surfboard.scenestateassertionservice</key>
		<true/>
		<key>com.apple.surfboard.sceneunderstandingservice</key>
		<true/>
		<key>com.apple.surfboard.searchservice</key>
		<true/>
		<key>com.apple.surfboard.sharingmodeservice</key>
		<true/>
		<key>com.apple.surfboard.shutdownservice</key>
		<true/>
		<key>com.apple.surfboard.testsupportservice</key>
		<true/>
		<key>com.apple.surfboard.virtualhandinteractionservice</key>
		<true/>
		<key>com.apple.systemstatus</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.systemstatus.activityattribution</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.systemstatus.publisher</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.usernotifications.remotenotificationservice</key>
		<true/>
		<key>com.apple.usernotifications.systemservice</key>
		<true/>
		<key>com.apple.usernotifications.usernotificationservice</key>
		<true/>
		<key>com.apple.usernotifications.usernotificationsettingsservice</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/SurfBoard.app/SurfBoard</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
	<key>_AdditionalProperties</key>
	<dict>
		<key>RunningBoard</key>
		<dict>
			<key>UnderlyingAssertion</key>
			<string>UnderlyingDarwinRoleUI</string>
		</dict>
	</dict>
</dict>
</plist>
