# beta/visionOS/System/Library/LaunchDaemons/com.apple.announced.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AppIntents</key>
	<dict>
		<key>MetadataAbsolutePaths</key>
		<array>
			<string>/System/Library/PrivateFrameworks/AnnounceDaemon.framework/</string>
		</array>
	</dict>
	<key>CFBundleIdentifier</key>
	<string>com.apple.announced</string>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.announced</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>com.apple.announce</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.announce</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.announce.activity.dailyanalytics</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>GracePeriod</key>
				<integer>300</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.announced.localplaybacksession</key>
		<true/>
		<key>com.apple.announced.reachability</key>
		<true/>
		<key>com.apple.announced.remoteplaybacksession</key>
		<true/>
		<key>com.apple.announced.server</key>
		<true/>
		<key>com.apple.announced.toneplayer</key>
		<true/>
		<key>com.apple.private.alloy.intercom-idswake</key>
		<true/>
		<key>com.apple.private.appintents.delegate.com.apple.announced</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/announced</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/announced</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
