# beta/visionOS/System/Library/LaunchDaemons/com.apple.pasteboard.pasted.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.pasteboard.pasted</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.coreservices.uauseractivitypasteboardclient.xpc</key>
		<true/>
		<key>com.apple.pasteboard.cache-delete</key>
		<true/>
		<key>com.apple.pasteboard.pasted</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/Pasteboard.framework/Support/pasted</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
