# beta/visionOS/System/Library/LaunchDaemons/com.apple.polarisd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<false/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.polarisd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.polaris.daemon_default</key>
		<true/>
		<key>com.apple.polaris.graphexecutor</key>
		<true/>
		<key>com.apple.polaris.graphstore</key>
		<true/>
		<key>com.apple.polaris.polarisd_debug</key>
		<true/>
		<key>com.apple.polaris.replay</key>
		<true/>
		<key>com.apple.polaris.shared_stream_manager</key>
		<true/>
		<key>com.apple.polaris.system</key>
		<true/>
		<key>com.apple.polaris.systemgraph</key>
		<true/>
		<key>com.apple.polaris.systemgraph_v2</key>
		<true/>
		<key>com.apple.polaris.systemtransitionservice</key>
		<true/>
		<key>com.apple.polaris.telemetry</key>
		<true/>
		<key>com.apple.polarissystemtransitiond.manifestagentservice</key>
		<true/>
	</dict>
	<key>OnDemand</key>
	<false/>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/Polaris.framework/polarisd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>_PanicOnCrash</key>
	<true/>
	<key>com.apple.private.spawn-panic-crash-behavior</key>
	<true/>
</dict>
</plist>
