# beta/visionOS/System/Library/LaunchDaemons/com.apple.AuthenticationServicesCore.AuthenticationServicesAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>EnabledPressureExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.AuthenticationServicesCore.AuthenticationServicesAgent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.AuthenticationServicesCore.AuthenticationServicesAgent.PeriodicMaintenanceActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>300</integer>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassA</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.AuthenticationServices.AuthenticationServicesAgent.CredentialUpdate</key>
		<true/>
		<key>com.apple.AuthenticationServices.AuthenticationServicesAgent.FileVaultRecoveryKeys</key>
		<true/>
		<key>com.apple.AuthenticationServices.AutoFill</key>
		<true/>
		<key>com.apple.AuthenticationServices.CredentialSharingGroups</key>
		<true/>
		<key>com.apple.AuthenticationServicesCore.AuthenticationServicesAgent</key>
		<true/>
		<key>com.apple.AuthenticationServicesCore.AuthenticationServicesAgent.CredentialExchange</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/AuthenticationServicesAgent</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
