# beta/visionOS/System/Library/LaunchDaemons/com.apple.DragUI.druid.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.DragUI.druid</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.DragUI.druid.destination</key>
		<true/>
		<key>com.apple.DragUI.druid.paste-announcement</key>
		<true/>
		<key>com.apple.DragUI.druid.source</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/DragUI.framework/Support/druid</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>3</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
