# beta/visionOS/System/Library/LaunchDaemons/com.apple.checkerboard.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>CFBundleIdentifier</key>
	<string>com.apple.CheckerBoard</string>
	<key>EmbeddedHomeScreen</key>
	<true/>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<false/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.checkerboard</string>
	<key>LimitLoadToHardware</key>
	<dict>
		<key>osenvironment</key>
		<array>
			<string>diagnostics</string>
		</array>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>PurpleSystemAppPort</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.CheckerBoard.services</key>
		<true/>
		<key>com.apple.UIKit.KeyboardManagement.hosted</key>
		<true/>
		<key>com.apple.UIKit.statusbarserver</key>
		<true/>
		<key>com.apple.frontboard.systemappservices</key>
		<true/>
		<key>com.apple.frontboard.workspace</key>
		<true/>
		<key>com.apple.private.corewifi.mobilewifi-xpc</key>
		<true/>
		<key>com.apple.systemstatus</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.systemstatus.activityattribution</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.systemstatus.publisher</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
	</dict>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/CheckerBoard.app/CheckerBoard</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>UserName</key>
	<string>mobile</string>
	<key>_AdditionalProperties</key>
	<dict>
		<key>RunningBoard</key>
		<dict>
			<key>UnderlyingAssertion</key>
			<string>UnderlyingDarwinRoleUI</string>
		</dict>
	</dict>
</dict>
</plist>
