# beta/visionOS/System/Library/LaunchDaemons/com.apple.jetpackassetd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.jetpackassetd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.jetpackassetd.maintenance.cleanup</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.jetpackassetd.task-group</string>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>604800</integer>
					<key>RandomInitialDelay</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
			</dict>
			<key>com.apple.jetpackassetd.maintenance.refresh</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.jetpackassetd.task-group</string>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>259200</integer>
					<key>RandomInitialDelay</key>
					<integer>14400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.jetpackassetd.postinstall</key>
			<dict>
				<key>NonRepeatingTask</key>
				<dict>
					<key>ScheduleAfter</key>
					<integer>60</integer>
					<key>TrySchedulingBefore</key>
					<integer>300</integer>
				</dict>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.jetpackassetd</key>
		<true/>
		<key>com.apple.jetpackassetd.xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/JetCore.framework/Support/jetpackassetd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
