# beta/visionOS/System/Library/LaunchDaemons/com.apple.roya.royacaptured.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>GroupName</key>
	<string>mobile</string>
	<key>Label</key>
	<string>com.apple.roya.royacaptured</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.roya.royacaptured</key>
		<false/>
		<key>com.apple.roya.royacaptured.camera</key>
		<false/>
		<key>com.apple.roya.royacaptured.system</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/royacaptured</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
