# beta/visionOS/System/Library/LaunchDaemons/com.apple.nfcd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>ExitTimeOut</key>
	<integer>12</integer>
	<key>Label</key>
	<string>com.apple.nfcd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>fury.spmi.activity</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>AppleStockholmSPMINotification</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.nfcd.coreanalytics.event</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>600</integer>
				<key>GracePeriod</key>
				<integer>300</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.nfcacd.hwmanager</key>
		<true/>
		<key>com.apple.nfcd.hwmanager</key>
		<true/>
		<key>com.apple.nfcd.service.corenfc</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/usr/libexec/nfcd</string>
	<key>PublishesEvents</key>
	<array>
		<string>com.apple.dockaccessoryd.pairsetup</string>
		<string>com.apple.nfcd.homekit.field</string>
		<string>com.apple.nfcd.homekit.proxcard</string>
		<string>com.apple.nfcd.xpc.sesd.events</string>
		<string>com.apple.accessoryd.nf-events</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
