# beta/visionOS/System/Library/LaunchDaemons/com.apple.realityenrolld.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.realityenrolld</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.realityenrolld.scheduledtasks.dailyanalytics</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>A</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.BiometricKit.enrollmentChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.BiometricKit.enrollmentChanged</string>
			</dict>
			<key>com.apple.Presence.identityProtectionLevelLowered</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.Presence.identityProtectionLevelLowered</string>
			</dict>
			<key>com.apple.managedconfiguration.passcodechanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.passcodechanged</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.realityenrolld.activity.postInstall</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.enroll.control</key>
		<true/>
		<key>com.apple.enroll.faceRepresentation</key>
		<true/>
		<key>com.apple.enroll.handSkinTone</key>
		<true/>
		<key>com.apple.enroll.handTexture</key>
		<true/>
		<key>com.apple.enroll.handshape</key>
		<true/>
		<key>com.apple.enroll.head</key>
		<true/>
		<key>com.apple.enroll.injection</key>
		<true/>
		<key>com.apple.enroll.state</key>
		<true/>
		<key>com.apple.enroll.userBody</key>
		<true/>
		<key>com.apple.enrollcore.deviceProtection</key>
		<true/>
		<key>com.apple.enrollcore.identity</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/realityenrolld</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
