# beta/visionOS/System/Library/LaunchDaemons/com.apple.finhealthd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.finhealthd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.finhealth.finance-background-delivery</key>
			<dict>
				<key>NonRepeatingTask</key>
				<dict>
					<key>ScheduleAfter</key>
					<integer>120</integer>
					<key>TrySchedulingBefore</key>
					<integer>7200</integer>
				</dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.finhealthd.finance.insights</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>69120</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.finhealthd.service</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/FinHealth.framework/finhealthd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
