# beta/visionOS/System/Library/LaunchDaemons/com.apple.airplayd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.apple.airplayd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.device-attach.b137</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IOUSBHostDevice</string>
				<key>idProduct</key>
				<integer>4781</integer>
				<key>idVendor</key>
				<integer>1452</integer>
			</dict>
			<key>com.apple.device-attach.b165</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IOUSBHostDevice</string>
				<key>idProduct</key>
				<integer>4782</integer>
				<key>idVendor</key>
				<integer>1452</integer>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.airplay.autoconnect.services</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.airplay.discoverybroker</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.airplay.endpoint.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.airplay.receiver.services</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.carkit.sessionRequestHandler</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.airplayd.audioprocessingtap.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.bufferedairplayglobalroutingregistry.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.carplayavvc.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.endpointmanager.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.endpointplaybacksession.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.endpointremotecontrolsession.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.endpointstream.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.endpointstreamaudioengine.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.nerotransportconnectionxpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.samplebufferconsumer.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/usr/libexec/airplayd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/airplayd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
