# beta/visionOS/System/Library/LaunchDaemons/com.apple.backboardd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<false/>
	<key>HighPriorityIO</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.backboardd</string>
	<key>MachServices</key>
	<dict>
		<key>PurpleSystemEventPort</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>PurpleWorkspacePort</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.AttentionAwareness</key>
		<true/>
		<key>com.apple.CARenderServer</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.accessibility.AXBackBoardServer</key>
		<true/>
		<key>com.apple.accessibility.gax.backboard</key>
		<true/>
		<key>com.apple.backboard.TouchDeliveryPolicyServer</key>
		<true/>
		<key>com.apple.backboard.altsysapp</key>
		<true/>
		<key>com.apple.backboard.display.services</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.backboard.hid-services.xpc</key>
		<true/>
		<key>com.apple.backboard.hid.services</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.backboard.oswatchdog</key>
		<true/>
		<key>com.apple.backboard.system-app-server</key>
		<true/>
		<key>com.apple.backboard.watchdog</key>
		<true/>
		<key>com.apple.backboardd.virtualframebuffer</key>
		<true/>
		<key>com.apple.backlightd</key>
		<true/>
		<key>com.apple.dmm.service</key>
		<true/>
		<key>com.apple.hid.PencilHaptics</key>
		<true/>
		<key>com.apple.iohideventsystem</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.iphone.axserver-systemwide</key>
		<true/>
		<key>com.apple.libEDR.Notification</key>
		<true/>
		<key>com.apple.lightsourcesupport.lightstate</key>
		<true/>
		<key>com.apple.realitysimulation.access</key>
		<true/>
		<key>com.apple.realitysimulation.animation-service</key>
		<true/>
		<key>com.apple.realitysimulation.apps</key>
		<true/>
		<key>com.apple.realitysimulation.bounds</key>
		<true/>
		<key>com.apple.realitysimulation.break</key>
		<true/>
		<key>com.apple.realitysimulation.connect</key>
		<true/>
		<key>com.apple.realitysimulation.data</key>
		<true/>
		<key>com.apple.realitysimulation.debug</key>
		<true/>
		<key>com.apple.realitysimulation.display</key>
		<true/>
		<key>com.apple.realitysimulation.effects</key>
		<true/>
		<key>com.apple.realitysimulation.host</key>
		<true/>
		<key>com.apple.realitysimulation.interaction</key>
		<true/>
		<key>com.apple.realitysimulation.naturalinput</key>
		<true/>
		<key>com.apple.realitysimulation.performance</key>
		<true/>
		<key>com.apple.realitysimulation.renderedcontentservice</key>
		<true/>
		<key>com.apple.realitysimulation.scene</key>
		<true/>
		<key>com.apple.realitysimulation.transitions</key>
		<true/>
		<key>com.apple.realitysimulation.vi</key>
		<true/>
		<key>com.apple.realitysystemsupport.hid_server_backboard</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/backboardd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
	<key>_Conclave</key>
	<string>com.apple.backboardd.conclave</string>
</dict>
</plist>
