# beta/visionOS/System/Library/LaunchDaemons/com.apple.parsec-fbf.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.parsec-fbf</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.parsec-fbf.cleanup</key>
			<dict>
				<key>ActivityGroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>ActivityGroupName</key>
				<string>com.apple.parsec-fbf.uploadActivities</string>
				<key>NetworkDownloadSize</key>
				<integer>524288</integer>
				<key>NetworkUploadSize</key>
				<integer>5242880</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>3600</integer>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>21600</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.parsec-fbf.flush</key>
			<dict>
				<key>ActivityGroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>ActivityGroupName</key>
				<string>com.apple.parsec-fbf.uploadActivities</string>
				<key>NetworkDownloadSize</key>
				<integer>524288</integer>
				<key>NetworkUploadSize</key>
				<integer>5242880</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>3600</integer>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>2700</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.LaunchServices.ApplicationsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LaunchServices.ApplicationsChanged</string>
			</dict>
			<key>com.apple.appprotection.change.hidden</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appprotection.change.hidden</string>
			</dict>
			<key>com.apple.parsec-fbf.FLUploadImmediately</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.parsec-fbf.FLUploadImmediately</string>
			</dict>
			<key>com.apple.siri.history.deletion.requested</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.history.deletion.requested</string>
			</dict>
			<key>com.apple.spotlightui.prefschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.spotlightui.prefschanged</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.feedbacklogger</key>
		<true/>
		<key>com.apple.parsec-fbf</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/CoreParsec.framework/parsec-fbf</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
