# beta/visionOS/System/Library/LaunchDaemons/com.apple.Safari.passwordbreachd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.Safari.passwordbreachd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.Safari.PasswordBreachHelper.PeriodicLookupActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>300</integer>
				<key>GracePeriod</key>
				<integer>172800</integer>
				<key>Interval</key>
				<integer>604800</integer>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassA</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.Safari.PasswordBreachHelper</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Background</string>
	<key>Program</key>
	<string>/usr/libexec/passwordbreachd</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
