# beta/visionOS/System/Library/LaunchDaemons/com.apple.email.maild.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.email.maild</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.alarm</key>
		<true/>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.email.categorization.FastPass</key>
			<dict>
				<key>FeatureCodes</key>
				<array>
					<integer>300</integer>
				</array>
				<key>NonRepeatingTask</key>
				<dict>
					<key>FastPassTask</key>
					<dict>
						<key>ProcessingTaskIdentifiers</key>
						<array>
							<string>com.apple.mail.blackpearl.categorizer</string>
						</array>
					</dict>
					<key>ScheduleAfter</key>
					<integer>300</integer>
					<key>TrySchedulingBefore</key>
					<integer>3600</integer>
				</dict>
				<key>RelatedApplications</key>
				<array>
					<string>com.apple.mobilemail</string>
					<string>com.apple.mail</string>
				</array>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RunWhenAppLaunchUnlikely</key>
				<true/>
			</dict>
			<key>com.apple.email.search.FastPass</key>
			<dict>
				<key>FeatureCodes</key>
				<array>
					<integer>301</integer>
				</array>
				<key>NonRepeatingTask</key>
				<dict>
					<key>FastPassTask</key>
					<dict>
						<key>ProcessingTaskIdentifiers</key>
						<array/>
					</dict>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.mail.blackpearl.categorizer</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RelatedApplications</key>
				<array>
					<string>com.apple.mobilemail</string>
					<string>com.apple.mail</string>
				</array>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>1800</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RunWhenAppLaunchUnlikely</key>
				<true/>
			</dict>
			<key>com.apple.mail.businessconnectgrouping</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RequiresSignificantUserInactivity</key>
				<false/>
			</dict>
			<key>com.apple.mail.fileProtectionClassMigration</key>
			<dict>
				<key>NonRepeatingTask</key>
				<dict>
					<key>TrySchedulingBefore</key>
					<integer>900</integer>
				</dict>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RelatedApplications</key>
				<array>
					<string>com.apple.mobilemail</string>
					<string>com.apple.mail</string>
				</array>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RunOnAppForeground</key>
				<false/>
			</dict>
			<key>com.apple.mail.messageauthentication</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RelatedApplications</key>
				<array>
					<string>com.apple.mobilemail</string>
					<string>com.apple.mail</string>
				</array>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>900</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RunOnAppForeground</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>MFNanoMailImportantBridgeSettingHasChangedDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>MFNanoMailImportantBridgeSettingHasChangedDarwinNotification</string>
			</dict>
			<key>PCPreferencesDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>PCPreferencesDidChangeNotification</string>
			</dict>
			<key>com.apple.dataaccess.ping</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.dataaccess.ping</string>
			</dict>
			<key>com.apple.exchangesyncd.ping</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.exchangesyncd.ping</string>
			</dict>
			<key>com.apple.mobilemail.afc.poll</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilemail.afc.poll</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.email.ProcessSQLQueryPerformanceData</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintainance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.mobilemail.activity.reset-device-activity-alarms</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>1</integer>
				<key>GracePeriod</key>
				<integer>1</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequiresClassA</key>
				<false/>
				<key>RequiresClassB</key>
				<false/>
				<key>RequiresClassC</key>
				<false/>
			</dict>
			<key>com.apple.mobilemail.scheduledFetches</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>600</integer>
				<key>Interval</key>
				<integer>900</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.corespotlight.daemon.mail</key>
		<true/>
		<key>com.apple.email.maild</key>
		<true/>
		<key>com.apple.email.maild.aps</key>
		<true/>
		<key>com.apple.mobilemail.nano.pairedsync</key>
		<true/>
		<key>com.apple.mobilemail.services.xpc</key>
		<true/>
		<key>com.apple.private.alloy.mail.fetches-idswake</key>
		<true/>
		<key>com.apple.private.alloy.mail.sync.accounts-idswake</key>
		<true/>
		<key>com.apple.private.alloy.mail.sync.content-idswake</key>
		<true/>
		<key>com.apple.private.alloy.mail.sync.messages-idswake</key>
		<true/>
		<key>com.apple.private.alloy.mail.sync.protected-idswake</key>
		<true/>
		<key>com.apple.private.alloy.mail.sync.protected.content-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.mobilemail</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/EmailDaemon.framework/maild</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
