# beta/visionOS/System/Library/LaunchDaemons/com.apple.corerepaird.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.corerepaird</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.corerepair</key>
		<true/>
		<key>com.apple.corerepair.diagnostics-controller</key>
		<true/>
		<key>com.apple.corerepair.intentControl</key>
		<true/>
		<key>com.apple.corerepair.preflightControl</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/corerepaird</string>
	</array>
	<key>_LimitLoadFromVariant</key>
	<string>IsDarwinOS</string>
</dict>
</plist>
