# beta/visionOS/System/Library/LaunchDaemons/com.apple.realitycamerad.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.realitycamerad</string>
	<key>LimitLoadToHardware</key>
	<dict>
		<key>targettype</key>
		<array>
			<string>N301</string>
			<string>N301a</string>
		</array>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.realitycamerad.RCPLSProvider</key>
		<true/>
		<key>com.apple.realitycamerad.RCXPCProvider</key>
		<true/>
		<key>com.apple.realitycamerad.watchdog</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/RealityCamera.framework/realitycamerad</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>_PanicOnCrash</key>
	<true/>
	<key>com.apple.private.spawn-panic-crash-behavior</key>
	<true/>
</dict>
</plist>
