# beta/visionOS/System/Library/LaunchDaemons/com.apple.fitnesscoachingd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>NSRunningFromLaunchd</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.fitnesscoachingd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.biome.compute.publisher</key>
		<dict>
			<key>com.apple.fitnesscoachingd.DeviceScreenLockedSubscription</key>
			<dict>
				<key>client</key>
				<string>com.apple.fitnesscoachingd</string>
				<key>graph</key>
				<dict>
					<key>dsl&lt;@BMDSLFilter&gt;</key>
					<dict>
						<key>comparison</key>
						<integer>1</integer>
						<key>keyPath</key>
						<string>self</string>
						<key>name</key>
						<string>filter</string>
						<key>negation</key>
						<false/>
						<key>upstream&lt;@BMDSLMap&gt;</key>
						<dict>
							<key>name</key>
							<string>map</string>
							<key>transform&lt;@BMDSLClassTransform&gt;</key>
							<dict>
								<key>className</key>
								<string>BMDSLTransforms</string>
								<key>name</key>
								<string>transformName</string>
								<key>transformName</key>
								<string>firstEventOfDayWithEvent:timestamp:lastTimestampBuffer:dayStartOffset:</string>
								<key>version</key>
								<integer>1</integer>
							</dict>
							<key>upstream&lt;@BMDSLMap&gt;</key>
							<dict>
								<key>name</key>
								<string>map</string>
								<key>transform&lt;@BMDSLArgumentTransform&gt;</key>
								<dict>
									<key>name</key>
									<string>argument</string>
									<key>selector&lt;@BMDSLKeyPathSelector&gt;</key>
									<dict>
										<key>keyPathSelector</key>
										<array>
											<string>self</string>
											<string>timestamp</string>
											<string>lastTimestampBuffer</string>
											<string>args.dayStartOffset</string>
										</array>
										<key>name</key>
										<string>keyPathSelector</string>
										<key>version</key>
										<integer>1</integer>
									</dict>
									<key>version</key>
									<integer>1</integer>
								</dict>
								<key>upstream&lt;@BMDSLFilter&gt;</key>
								<dict>
									<key>comparison</key>
									<integer>0</integer>
									<key>keyPath</key>
									<string>eventBody.starting</string>
									<key>name</key>
									<string>filter</string>
									<key>negation</key>
									<false/>
									<key>upstream&lt;@BMDSLState&gt;</key>
									<dict>
										<key>key</key>
										<string>args</string>
										<key>name</key>
										<string>state</string>
										<key>upstream&lt;@BMDSLState&gt;</key>
										<dict>
											<key>key</key>
											<string>lastTimestampBuffer</string>
											<key>name</key>
											<string>state</string>
											<key>upstream&lt;@BMDSLStreamPublisher&gt;</key>
											<dict>
												<key>name</key>
												<string>streamPublisher</string>
												<key>streamIdentifier</key>
												<string>Device.ScreenLocked</string>
												<key>streamType</key>
												<integer>2</integer>
												<key>version</key>
												<integer>1</integer>
											</dict>
											<key>value</key>
											<dict/>
											<key>version</key>
											<integer>1</integer>
										</dict>
										<key>value</key>
										<dict>
											<key>dayStartOffset</key>
											<integer>10800</integer>
										</dict>
										<key>version</key>
										<integer>1</integer>
									</dict>
									<key>value</key>
									<false/>
									<key>version</key>
									<integer>1</integer>
								</dict>
								<key>version</key>
								<integer>1</integer>
							</dict>
							<key>version</key>
							<integer>1</integer>
						</dict>
						<key>value</key>
						<false/>
						<key>version</key>
						<integer>1</integer>
					</dict>
				</dict>
				<key>identifier</key>
				<string>com.apple.fitnesscoachingd.DeviceScreenLockedSubscription</string>
			</dict>
		</dict>
		<key>com.apple.coreduetcontext.client_event_stream</key>
		<dict>
			<key>notification</key>
			<dict>
				<key>client-identifier</key>
				<string>com.apple.fitnesscoachingd.ContextStoreClient</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>FitnessPlusPlanCoachingDefaultsUpdatedNotification</key>
			<dict>
				<key>Notification</key>
				<string>FitnessPlusPlanCoachingDefaultsUpdatedNotification</string>
			</dict>
			<key>NanoLifestylePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>NanoLifestylePreferencesChangedNotification</string>
			</dict>
			<key>SeymourWorkoutPlanChanged</key>
			<dict>
				<key>Notification</key>
				<string>SeymourWorkoutPlanChanged</string>
			</dict>
			<key>SignificantTimeChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>SignificantTimeChangeNotification</string>
			</dict>
			<key>com.apple.FitnessCoaching.FCPauseRingsSampleChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.FitnessCoaching.FCPauseRingsSampleChangedNotification</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
			<key>com.apple.nanoregistry.watchdidbecomeactive</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.watchdidbecomeactive</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.fitnesscoaching.Scheduling.InitialLaunchActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>300</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<false/>
				<key>RequiresClassC</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.fitnesscoachingd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/FitnessCoachingServices.framework/fitnesscoachingd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
