# beta/visionOS/System/Library/LaunchDaemons/com.apple.dt.previewsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.previewsd</string>
	<key>LaunchEvents</key>
	<dict/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.dt.previewsd.agentpipeservice</key>
		<true/>
		<key>com.apple.dt.previewsd.agentpipeservice.host-to-sim</key>
		<true/>
		<key>com.apple.dt.previewsd.crashlistener</key>
		<true/>
		<key>com.apple.dt.previewsd.installservice</key>
		<true/>
		<key>com.apple.dt.previewsd.previewservice</key>
		<true/>
		<key>com.apple.dt.previewsd.previewservice.host-to-sim</key>
		<true/>
		<key>com.apple.dt.previewsd.processexit</key>
		<true/>
		<key>com.apple.dt.previewsd.shellpipeservice</key>
		<true/>
		<key>com.apple.dt.previewsd.shellpipeservice.host-to-sim</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/PreviewsOSSupport.framework/Support/previewsd</string>
	<key>ProgramArguments</key>
	<array/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
