# diff/iOS/System/Library/LaunchAngels/com.apple.GameOverlayUI.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<false/>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.apple.GameOverlayUI</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.GameOverlayUI</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
	</dict>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/GameOverlayUI.app/GameOverlayUI</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
	<key>_AdditionalProperties</key>
	<dict>
		<key>RunningBoard</key>
		<dict>
			<key>Managed</key>
			<true/>
			<key>Reported</key>
			<true/>
		</dict>
	</dict>
</dict>
</plist>
