# diff/iOS/System/Library/LaunchDaemons/com.apple.mobile_installation_proxy.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.mobile.installation_proxy</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mobile.installation_proxy</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/mobile_installation_proxy</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
