# diff/iOS/System/Library/LaunchDaemons/com.apple.nanotimekitcompaniond.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.nanotimekitcompaniond</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>LanguageChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.language.changed</string>
			</dict>
			<key>LibraryCollectionTargetMapDataChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanophotos.prefs.LibraryCollectionTargetMapData-changed</string>
			</dict>
			<key>LocaleChange</key>
			<dict>
				<key>Notification</key>
				<string>AppleDatePreferencesChangedNotification</string>
			</dict>
			<key>NanoPhotosLibraryChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.NanoPhotos.Library.changed</string>
			</dict>
			<key>com.apple.LaunchServices.ApplicationsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LaunchServices.ApplicationsChanged</string>
			</dict>
			<key>devices-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.replicatord.devicesChanged</string>
			</dict>
			<key>messages-received</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanotimekit.replicator.library.message</string>
			</dict>
			<key>records-received</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanotimekit.replicator.library.record</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.nanotimekit.facesupport</key>
		<true/>
		<key>com.apple.clockface.pairedsync</key>
		<true/>
		<key>com.apple.nanofacegallery.allassertions.server</key>
		<true/>
		<key>com.apple.nanofacegallery.gallerybroker.server</key>
		<true/>
		<key>com.apple.nanofacegallery.librarybroker.server</key>
		<true/>
		<key>com.apple.nanotimekit.collectionserver</key>
		<true/>
		<key>com.apple.nanotimekit.complicationstoreserver</key>
		<true/>
		<key>com.apple.nanotimekit.devicelibraryserver</key>
		<true/>
		<key>com.apple.nanotimekit.facesnapshotserver</key>
		<true/>
		<key>com.apple.nanotimekit.facesupportserver</key>
		<true/>
		<key>com.apple.nanotimekit.librarysnapshotreplicator</key>
		<true/>
		<key>com.apple.nanotimekit.replication.library.listener</key>
		<true/>
		<key>com.apple.nanotimekit.trackedphotosserver</key>
		<true/>
		<key>com.apple.private.alloy.clockface.sync-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/NanoTimeKit.framework/nanotimekitcompaniond</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
