# diff/iOS/System/Library/LaunchDaemons/com.apple.coreidvd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.coreidvd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.accounts.changes</key>
		<dict>
			<key>AppleAccountAddedOrDeleted</key>
			<dict>
				<key>com.apple.account.AppleAccount</key>
				<dict>
					<key>ChangeTypes</key>
					<array>
						<integer>1</integer>
						<integer>3</integer>
					</array>
				</dict>
			</dict>
		</dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.coreidvd.biome</key>
			<dict>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>TrySchedulingBefore</key>
					<integer>28800</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
			</dict>
			<key>com.apple.coreidvd.dpdt</key>
			<dict>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>604800</integer>
					<key>TrySchedulingBefore</key>
					<integer>28800</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
			</dict>
			<key>com.apple.coreidvd.gc</key>
			<dict>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>TrySchedulingBefore</key>
					<integer>28800</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
			</dict>
			<key>com.apple.coreidvd.gcmdr</key>
			<dict>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>TrySchedulingBefore</key>
					<integer>28800</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
			</dict>
			<key>com.apple.coreidvd.idp</key>
			<dict>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>TrySchedulingBefore</key>
					<integer>28800</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.coreidvd.wu</key>
			<dict>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>TrySchedulingBefore</key>
					<integer>28800</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.coreidvd.wudu</key>
			<dict>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>TrySchedulingBefore</key>
					<integer>28800</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>A</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>Application Installed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>Application Uninstalled</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.passd.matching</key>
		<dict>
			<key>Identity Provisioning Complete</key>
			<dict>
				<key>Name</key>
				<string>PKPassLibraryPendingProvisioningFinishedNotification</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>com.apple.coreidv.coreidvd.handoff</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.coreidv.coreidvd.handoff</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.coreidvd</key>
		<true/>
		<key>com.apple.coreidvd.digital-presentment.xpc</key>
		<true/>
		<key>com.apple.coreidvd.docUpload.xpc</key>
		<true/>
		<key>com.apple.coreidvd.identity-proofing-data-sharing.xpc</key>
		<true/>
		<key>com.apple.coreidvd.identity-proofing-ui.xpc</key>
		<true/>
		<key>com.apple.coreidvd.identity-proofing.xpc</key>
		<true/>
		<key>com.apple.coreidvd.identity-provisioning.xpc</key>
		<true/>
		<key>com.apple.coreidvd.mobile-document-provider-registration.xpc</key>
		<true/>
		<key>com.apple.coreidvd.mobile-document-reader.xpc</key>
		<true/>
		<key>com.apple.coreidvd.proofing</key>
		<true/>
		<key>com.apple.coreidvd.report-a-concern.xpc</key>
		<true/>
		<key>com.apple.coreidvd.watch-notifications</key>
		<true/>
		<key>com.apple.coreidvd.web-presentment-remote-alert.xpc</key>
		<true/>
		<key>com.apple.coreidvd.web-presentment.xpc</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.coreidv.CoreIDVUserNotifications</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/coreidvd</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
