# diff/iOS/System/Library/LaunchDaemons/com.apple.filesystems.userfsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.filesystems.userfsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.userfsd.ahci.device.attached</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOParentMatch</key>
				<dict>
					<key>IOPropertyMatch</key>
					<dict>
						<key>Protocol Characteristics</key>
						<dict>
							<key>Physical Interconnect</key>
							<string>SATA</string>
							<key>Physical Interconnect Location</key>
							<string>External</string>
						</dict>
					</dict>
					<key>IOProviderClass</key>
					<string>IOBlockStorageDevice</string>
				</dict>
				<key>IOProviderClass</key>
				<string>IOMedia</string>
			</dict>
			<key>com.apple.userfsd.pcie.device.attached</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOParentMatch</key>
				<dict>
					<key>IOPropertyMatch</key>
					<dict>
						<key>Protocol Characteristics</key>
						<dict>
							<key>Physical Interconnect</key>
							<string>PCI-Express</string>
							<key>Physical Interconnect Location</key>
							<string>External</string>
						</dict>
					</dict>
					<key>IOProviderClass</key>
					<string>IOBlockStorageDevice</string>
				</dict>
				<key>IOProviderClass</key>
				<string>IOMedia</string>
			</dict>
			<key>com.apple.userfsd.usb.device.attached</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOParentMatch</key>
				<dict>
					<key>IOPropertyMatch</key>
					<dict>
						<key>removable</key>
						<true/>
					</dict>
					<key>IOProviderClass</key>
					<string>IOUSBHostDevice</string>
				</dict>
				<key>IOProviderClass</key>
				<string>IOBlockStorageDevice</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.filesystems.localLiveFiles</key>
		<true/>
		<key>com.apple.filesystems.userfsd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/UserFS.framework/userfsd</string>
		<string>-c</string>
		<string>-S</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>WorkingDirectory</key>
	<string>/System/Library/PrivateFrameworks/UserFS.framework/PlugIns</string>
</dict>
</plist>
