# diff/iOS/System/Library/LaunchDaemons/com.apple.healthd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>CFBundleAllowMixedLocalizations</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.healthd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.alarm</key>
		<true/>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.healthd.post-install-update.fastpass</key>
			<dict>
				<key>FeatureCodes</key>
				<array>
					<integer>600</integer>
				</array>
				<key>NonRepeatingTask</key>
				<dict>
					<key>FastPassTask</key>
					<dict/>
				</dict>
				<key>PostInstall</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>B</string>
			</dict>
			<key>com.apple.healthd.sync.cloud.initial-restore.fastpass</key>
			<dict>
				<key>Dependencies</key>
				<array>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.healthd.post-install-update.fastpass</string>
						<key>DependencyType</key>
						<string>DependencyTypeCompletion</string>
					</dict>
				</array>
				<key>FeatureCodes</key>
				<array>
					<integer>601</integer>
				</array>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.healthd.sync.cloud</string>
				<key>NonRepeatingTask</key>
				<dict>
					<key>FastPassTask</key>
					<dict>
						<key>ProcessingTaskIdentifiers</key>
						<array>
							<string>com.apple.healthd.sync.cloud.periodic</string>
						</array>
					</dict>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>B</string>
			</dict>
			<key>com.apple.healthd.sync.cloud.periodic</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.healthd.sync.cloud</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>14400</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>B</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>ConnectedGymPreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>ConnectedGymPreferencesChangedNotification</string>
			</dict>
			<key>com.apple.Preferences.ResetPrivacyWarningsNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.Preferences.ResetPrivacyWarningsNotification</string>
			</dict>
			<key>com.apple.bluetooth.accessory-authentication.success</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.bluetooth.accessory-authentication.success</string>
			</dict>
			<key>com.apple.nearfield.handoff.terminal</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nearfield.handoff.terminal</string>
			</dict>
			<key>com.apple.siri.koa.donate</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.koa.donate</string>
			</dict>
			<key>com.apple.system.timezone</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.timezone</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>HKRapportEvent</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.healthd.rapport</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.healthd.prune-database</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>14400</integer>
				<key>Interval</key>
				<integer>28800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RequiresClassA</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.ActivitySharing.push</key>
		<true/>
		<key>com.apple.aps.health.records.push</key>
		<true/>
		<key>com.apple.aps.healthd.cloud.subscription.push</key>
		<true/>
		<key>com.apple.healthd.cache-delete</key>
		<true/>
		<key>com.apple.healthd.server</key>
		<true/>
		<key>com.apple.private.alloy.health.sync.classc-idswake</key>
		<true/>
		<key>com.apple.private.alloy.healthapp.messaging</key>
		<true/>
		<key>com.apple.private.alloy.healthappnotificationsync-idswake</key>
		<true/>
		<key>com.apple.private.healthd.pairedsync</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/HealthKit.framework/healthd</string>
	</array>
	<key>PublishesEvents</key>
	<string>com.apple.healthd.authorization_events</string>
	<key>RunAtLoad</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>0</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
