# diff/iOS/System/Library/LaunchDaemons/com.apple.ciphermld.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.ciphermld</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.CipherML.CacheCleaner</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>RandomInitialDelay</key>
					<integer>3600</integer>
				</dict>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.CipherML.RotationTask</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>43200</integer>
					<key>RandomInitialDelay</key>
					<integer>3600</integer>
				</dict>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.ciphermld</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/ciphermld</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/ciphermld</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
