# diff/iOS/System/Library/LaunchDaemons/com.apple.ap.promotedcontentd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.ap.promotedcontentd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.ap.promotedcontentd.jetpack.cache-daily</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.ap.promotedcontentd.jetpack.cache-group</string>
				<key>InvolvedProcesses</key>
				<array>
					<string>com.apple.ap.PromotedContentJetService</string>
				</array>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<string>86400</string>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.ap.promotedcontentd.postinstall</key>
			<dict>
				<key>NonRepeatingTask</key>
				<dict>
					<key>ScheduleAfter</key>
					<string>0</string>
				</dict>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.ap.adprivacyd.iTunesActiveStorefrontDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ap.adprivacyd.iTunesActiveStorefrontDidChangeNotification</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.ap.promotedcontentd.gc</key>
			<dict>
				<key>DiskIntensive</key>
				<true/>
				<key>GracePeriod</key>
				<integer>86000</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<false/>
				<key>RequireScreenSleep</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
			</dict>
			<key>com.apple.ap.promotedcontentd.timespent</key>
			<dict>
				<key>DiskIntensive</key>
				<true/>
				<key>GracePeriod</key>
				<integer>43200</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<integer>3600</integer>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.ap.adservicesd.statusconditionservice</key>
		<true/>
		<key>com.apple.ap.promotedcontent.attributionservice</key>
		<true/>
		<key>com.apple.ap.promotedcontent.identifiermanager</key>
		<true/>
		<key>com.apple.ap.promotedcontent.mescalinterface</key>
		<true/>
		<key>com.apple.ap.promotedcontent.metrics</key>
		<true/>
		<key>com.apple.ap.promotedcontent.pcinterface</key>
		<true/>
		<key>com.apple.ap.promotedcontent.signing</key>
		<true/>
		<key>com.apple.ap.promotedcontent.supportinterface</key>
		<true/>
		<key>com.apple.ap.promotedcontent.trialbucketid</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/promotedcontentd</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
