# diff/iOS/System/Library/LaunchDaemons/com.apple.sirittsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.sirittsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.sirittsd.cleanUp</key>
			<dict>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.sirittsd.neuralCompiling</key>
			<dict>
				<key>Interval</key>
				<integer>14400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireSignificantUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.sirittsd.postInstall</key>
			<dict>
				<key>PostInstall</key>
				<true/>
			</dict>
			<key>com.apple.sirittsd.voiceUpdate</key>
			<dict>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<integer>3600</integer>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.sirittsd.weeklyVoiceUpdate</key>
			<dict>
				<key>Interval</key>
				<integer>604800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.sirittsd</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/SiriTTSService.framework/sirittsd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
