# diff/iOS/System/Library/LaunchDaemons/com.apple.sharingd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.sharingd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>com.apple.DDUI-Picker</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.DDUI-Picker</string>
				<key>type</key>
				<string>server</string>
			</dict>
			<key>com.apple.devicediscoveryui.rapportwake</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.devicediscoveryui.rapportwake</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.DDUI.UserAcceptService</key>
		<true/>
		<key>com.apple.DeviceDiscoveryUI.Permission</key>
		<true/>
		<key>com.apple.SharingServices</key>
		<true/>
		<key>com.apple.coreservices.appleid.aps</key>
		<true/>
		<key>com.apple.private.alloy.continuity.auth-idswake</key>
		<true/>
		<key>com.apple.private.alloy.continuity.unlock-idswake</key>
		<true/>
		<key>com.apple.private.alloy.sharing.ranging-idswake</key>
		<true/>
		<key>com.apple.sharing.HUD</key>
		<true/>
		<key>com.apple.sharing.accessibility</key>
		<true/>
		<key>com.apple.sharing.adnotice.services</key>
		<true/>
		<key>com.apple.sharing.airdrop.service</key>
		<true/>
		<key>com.apple.sharing.classroom-transfer</key>
		<true/>
		<key>com.apple.sharing.collaborationuserdefaults</key>
		<true/>
		<key>com.apple.sharing.handoff.advertising</key>
		<true/>
		<key>com.apple.sharing.handoff.scanning</key>
		<true/>
		<key>com.apple.sharing.pinpairing</key>
		<true/>
		<key>com.apple.sharing.remote-text-editing</key>
		<true/>
		<key>com.apple.sharing.sharesheet</key>
		<true/>
		<key>com.apple.sharing.sharesheetrecipients</key>
		<true/>
		<key>com.apple.sharing.sharesheetuserdefaults</key>
		<true/>
		<key>com.apple.sharing.transfer-observer</key>
		<true/>
		<key>com.apple.sharingd</key>
		<true/>
		<key>com.apple.sharingd.b332setup</key>
		<true/>
		<key>com.apple.sharingd.nsxpc</key>
		<true/>
		<key>com.apple.sharingd.pairedcontactmanager</key>
		<true/>
		<key>com.apple.sharingd.paireddevice</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.SFUserNotificationTVLatency</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.Sharing.HomePodRemoteNotifications</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.Sharing.TVRemoteNotifications</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.Sharing.VisionUnlockiOSNotifications</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.Sharing.WatchRemoteNotifications</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.Sharing.b7ce1794a1c8766816fc7b7500742862</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.UserNotifications.DeviceDiscoveryUIAgent</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/usr/libexec/sharingd</string>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
