# diff/iOS/System/Library/LaunchDaemons/com.apple.PerformanceTrace.ptpassivecollectiond.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.PerformanceTrace.ptpassivecollectiond</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.performancetrace.metric_monitoring_timed_out</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.perfpowermetricmonitor.monitoring_timed_out</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.PerformanceTrace.ptpassivecollectiond.collectionconfig</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/ptpassivecollectiond</string>
	</array>
</dict>
</plist>
