# diff/iOS/System/Library/LaunchDaemons/com.apple.powerexperienced.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.powerexperienced</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.coreaudio.IORunning</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.coreaudio.IORunning</string>
			</dict>
			<key>com.apple.osintelligence.iblm.mitigationchanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.osintelligence.iblm.mitigationchanged</string>
			</dict>
			<key>com.apple.powerexperienced.testdevicecontext.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.powerexperienced.testdevicecontext.changed</string>
			</dict>
			<key>com.apple.powerexperienced.testmobilechargingcontroller.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.powerexperienced.testmobilechargingcontroller.changed</string>
			</dict>
			<key>com.apple.request.hipuncap</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.request.hipuncap</string>
			</dict>
			<key>com.apple.springboard.lockstate</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.lockstate</string>
			</dict>
			<key>com.apple.system.earlythermalnotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.earlythermalnotification</string>
			</dict>
			<key>com.apple.system.lowpowermode</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.lowpowermode</string>
			</dict>
			<key>com.apple.system.powermanagement.clamshellstate</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powermanagement.clamshellstate</string>
			</dict>
			<key>com.apple.system.powersources.source</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powersources.source</string>
			</dict>
			<key>com.apple.system.powersources.timeremaining</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powersources.timeremaining</string>
			</dict>
			<key>com.apple.system.thermalpressurelevel</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.thermalpressurelevel</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.powerexperienced.powermitigationsmanager.service</key>
		<true/>
		<key>com.apple.powerexperienced.powermodes</key>
		<true/>
		<key>com.apple.powerexperienced.powermodesservice</key>
		<true/>
		<key>com.apple.powerexperienced.powertargetcontroller</key>
		<true/>
		<key>com.apple.powerexperienced.powertargetcontrollerservice</key>
		<true/>
		<key>com.apple.powerexperienced.resourceusage</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/usr/libexec/powerexperienced</string>
</dict>
</plist>
