# diff/iOS/System/Library/LaunchDaemons/com.apple.biomed.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.biomed</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.cascade.database-maintenance.nightly-task</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AFAssistantEnablementDidChangeDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>AFAssistantEnablementDidChangeDarwinNotification</string>
			</dict>
			<key>com.apple.Preferences.ResetPrivacyWarningsNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.Preferences.ResetPrivacyWarningsNotification</string>
			</dict>
			<key>com.apple.TextInput.RemoveAllDynamicDictionariesNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TextInput.RemoveAllDynamicDictionariesNotification</string>
			</dict>
			<key>com.apple.siri.history.deletion.requested</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.history.deletion.requested</string>
			</dict>
			<key>com.apple.suggestions.settingsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.suggestions.settingsChanged</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.biome.prune-expired-events</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.ResetDisabledViewsDaily</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>AllowBattery</key>
				<false/>
				<key>CPUIntensive</key>
				<false/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.ViewDaily</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>CPUIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.ViewDailyWithNetworking</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>CPUIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.ViewEvery2Hours</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>AllowBattery</key>
				<true/>
				<key>CPUIntensive</key>
				<false/>
				<key>Interval</key>
				<integer>7200</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.ViewEvery4Hours</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>AllowBattery</key>
				<false/>
				<key>CPUIntensive</key>
				<false/>
				<key>Interval</key>
				<integer>14400</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.ViewWeekly</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>CPUIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>604800</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.biome.PublicStreamAccessService</key>
		<true/>
		<key>com.apple.biome.access.system</key>
		<true/>
		<key>com.apple.biome.access.user</key>
		<true/>
		<key>com.apple.biome.compute.publisher.service</key>
		<true/>
		<key>com.apple.biome.compute.publisher.service.user</key>
		<true/>
		<key>com.apple.biome.compute.source</key>
		<true/>
		<key>com.apple.biome.compute.source.user</key>
		<true/>
		<key>com.apple.cascade.Maintenance</key>
		<true/>
		<key>com.apple.cascade.SetChangeRelayService</key>
		<true/>
		<key>com.apple.intelligenceplatform.Coordination</key>
		<true/>
		<key>com.apple.intelligenceplatform.InternalBiome</key>
		<true/>
		<key>com.apple.intelligenceplatform.View</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/BiomeStreams.framework/Support/biomed</string>
	</array>
	<key>PublishesEvents</key>
	<dict>
		<key>com.apple.biome.compute.publisher</key>
		<dict>
			<key>RequireTrusted</key>
			<true/>
		</dict>
		<key>com.apple.biome.compute.publisher.user</key>
		<dict>
			<key>RequireTrusted</key>
			<true/>
		</dict>
		<key>com.apple.cascade.xpc_event.setChange</key>
		<dict>
			<key>RequireTrusted</key>
			<true/>
		</dict>
	</dict>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
