# diff/iOS/System/Library/LaunchDaemons/com.apple.icloud.findmydeviced.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.icloud.findmydeviced</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AppleDatePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleDatePreferencesChangedNotification</string>
			</dict>
			<key>AppleKeyboardsPreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleKeyboardsPreferencesChangedNotification</string>
			</dict>
			<key>AppleLanguagePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleLanguagePreferencesChangedNotification</string>
			</dict>
			<key>AppleNumberPreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleNumberPreferencesChangedNotification</string>
			</dict>
			<key>AppleTimePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleTimePreferencesChangedNotification</string>
			</dict>
			<key>com.apple.DuetHeuristic-BM.shutdowsoon</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.DuetHeuristic-BM.shutdowsoon</string>
			</dict>
			<key>com.apple.accessories.connection.MFi4AccessoryDisconnected</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessories.connection.MFi4AccessoryDisconnected</string>
			</dict>
			<key>com.apple.accessories.connection.passedMFi4Auth</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessories.connection.passedMFi4Auth</string>
			</dict>
			<key>com.apple.bluetooth.connection</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.bluetooth.connection</string>
			</dict>
			<key>com.apple.bluetooth.pairing</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.bluetooth.pairing</string>
			</dict>
			<key>com.apple.icloud.FindMy.addMagSafeAccessory</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.FindMy.addMagSafeAccessory</string>
			</dict>
			<key>com.apple.icloud.findmydeviced.findkit.magSafe.added</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.findmydeviced.findkit.magSafe.added</string>
			</dict>
			<key>com.apple.icloud.findmydeviced.findkit.magSafe.attach</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.findmydeviced.findkit.magSafe.attach</string>
			</dict>
			<key>com.apple.icloud.findmydeviced.findkit.magSafe.detach</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.findmydeviced.findkit.magSafe.detach</string>
			</dict>
			<key>com.apple.icloud.findmydeviced.findkit.magSafe.removed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.findmydeviced.findkit.magSafe.removed</string>
			</dict>
			<key>com.apple.icloud.searchparty.accessoryDidPair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.searchparty.accessoryDidPair</string>
			</dict>
			<key>com.apple.icloud.searchparty.secureLocations.liteLocationPublishRequest</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.searchparty.secureLocations.liteLocationPublishRequest</string>
			</dict>
			<key>com.apple.icloud.searchparty.selfbeaconchanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.searchparty.selfbeaconchanged</string>
			</dict>
			<key>com.apple.language.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.language.changed</string>
			</dict>
			<key>com.apple.locationd.authorization</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.authorization</string>
			</dict>
			<key>com.apple.locationd/Prefs</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd/Prefs</string>
			</dict>
			<key>com.apple.managedconfiguration.allowpasscodemodificationchanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.allowpasscodemodificationchanged</string>
			</dict>
			<key>com.apple.managedconfiguration.passcodechanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.passcodechanged</string>
			</dict>
			<key>com.apple.mobile.lockdown.activation_state</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.lockdown.activation_state</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidunpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidunpair</string>
			</dict>
			<key>com.apple.newdeviceoutreach.coverageupdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.newdeviceoutreach.coverageupdated</string>
			</dict>
			<key>com.apple.purplebuddy.setupdone</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.purplebuddy.setupdone</string>
			</dict>
			<key>com.apple.system.hostname</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.hostname</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.icloud.findmydeviced.launch-on-reboot</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>120</integer>
				<key>GracePeriod</key>
				<integer>30</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.icloud.findmydeviced</key>
		<true/>
		<key>com.apple.icloud.findmydeviced.aps-demo</key>
		<true/>
		<key>com.apple.icloud.findmydeviced.aps-development</key>
		<true/>
		<key>com.apple.icloud.findmydeviced.aps-production</key>
		<true/>
		<key>com.apple.private.alloy.findmydeviced.watch-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/findmydeviced</string>
	<key>Umask</key>
	<integer>23</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
