# diff/iOS/System/Library/LaunchDaemons/com.apple.homeenergyd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.homeenergyd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.alarm</key>
		<true/>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.homeenergyd.sync.activity</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>46800</integer>
				</dict>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.homeenergyd.firstlaunch</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<integer>1800</integer>
				<key>Repeating</key>
				<false/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<false/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.homeenergyd</key>
		<true/>
		<key>com.apple.energykit.localStore</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>HomeEnergy/xpcStore2025</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.energykit.privateStore</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>HomeEnergy/xpcStore2025</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.energykit.sharedStore</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>HomeEnergy/xpcStore2025</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.energykit.store</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>HomeEnergy/xpcStore2025</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.homeenergyd.xpc</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/HomeEnergyDaemon.framework/Support/homeenergyd</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
