# diff/iOS/System/Library/LaunchDaemons/com.apple.iosdiagnosticsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.iosdiagnosticsd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.iosdiagnosticsd</key>
		<true/>
		<key>com.apple.private.alloy.iosdiagnostics-idswake</key>
		<true/>
		<key>com.apple.private.alloy.iosdiagnostics.icloud-idswake</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/iOSDiagnostics.framework/iosdiagnosticsd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
