# diff/iOS/System/Library/LaunchDaemons/com.apple.ospredictiond.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.ospredictiond</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.coreduetcontext.client_event_stream</key>
		<dict>
			<key>notification</key>
			<dict>
				<key>client-identifier</key>
				<string>com.apple.osintelligence.iblm.contextstore-registration</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.system.console_mode_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.console_mode_changed</string>
			</dict>
			<key>kPLTaskingStartNotificationGlobal</key>
			<dict>
				<key>Notification</key>
				<string>kPLTaskingStartNotificationGlobal</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.osintelligence.evaluateModelType</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.osintelligence.inactivitybackup</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>345600</integer>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.OSIntelligence</key>
		<true/>
		<key>com.apple.OSIntelligence.battery</key>
		<true/>
		<key>com.apple.OSIntelligence.charging</key>
		<true/>
		<key>com.apple.OSIntelligence.lastlock</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.osintelligence.notifications</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/ospredictiond</string>
</dict>
</plist>
