# diff/iOS/System/Library/LaunchDaemons/com.apple.appstored.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.appstored</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
			<key>com.apple.ReportCrash.beta.bulletin</key>
			<dict>
				<key>Name</key>
				<string>com.apple.ReportCrash.beta.bulletin</string>
			</dict>
			<key>com.apple.StorageSettings.AppMayUninstall</key>
			<dict>
				<key>Name</key>
				<string>com.apple.StorageSettings.AppMayUninstall</string>
			</dict>
			<key>com.apple.appstored.ASDAppleIDDidChangeNotification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.appstored.ASDAppleIDDidChangeNotification</string>
			</dict>
			<key>com.apple.appstored.ASDItunesAccountDidChangeNotification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.appstored.ASDItunesAccountDidChangeNotification</string>
			</dict>
			<key>com.apple.appstored.ASDSandboxAccountDidChangeNotification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.appstored.ASDSandboxAccountDidChangeNotification</string>
			</dict>
			<key>com.apple.installcoordination.AppMayUninstall</key>
			<dict>
				<key>Name</key>
				<string>com.apple.installcoordination.AppMayUninstall</string>
			</dict>
			<key>com.apple.managedconfiguration.effectivesettingschanged</key>
			<dict>
				<key>Name</key>
				<string>com.apple.managedconfiguration.effectivesettingschanged</string>
			</dict>
			<key>com.lsinstallprogress.appcontrols.cancel</key>
			<dict>
				<key>Name</key>
				<string>com.lsinstallprogress.appcontrols.cancel</string>
			</dict>
			<key>com.lsinstallprogress.appcontrols.pause</key>
			<dict>
				<key>Name</key>
				<string>com.lsinstallprogress.appcontrols.pause</string>
			</dict>
			<key>com.lsinstallprogress.appcontrols.prioritize</key>
			<dict>
				<key>Name</key>
				<string>com.lsinstallprogress.appcontrols.prioritize</string>
			</dict>
			<key>com.lsinstallprogress.appcontrols.resume</key>
			<dict>
				<key>Name</key>
				<string>com.lsinstallprogress.appcontrols.resume</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.Preferences.ResetPrivacyWarningsNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.Preferences.ResetPrivacyWarningsNotification</string>
			</dict>
			<key>com.apple.StoreServices.StorefrontChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.StoreServices.StorefrontChanged</string>
			</dict>
			<key>com.apple.datamigrator.migrationDidFinish</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.datamigrator.migrationDidFinish</string>
			</dict>
			<key>com.apple.fairplayd.resync-fpkeybag</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.fairplayd.resync-fpkeybag</string>
			</dict>
			<key>com.apple.itunesstored.invalidatebags</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.itunesstored.invalidatebags</string>
			</dict>
			<key>com.apple.private.restrict-post.MobileBackup.backgroundCellularAccessChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.private.restrict-post.MobileBackup.backgroundCellularAccessChanged</string>
			</dict>
			<key>com.apple.purplebuddy.setupdone</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.purplebuddy.setupdone</string>
			</dict>
		</dict>
		<key>com.apple.nsurlsessiond-events</key>
		<dict>
			<key>LaunchOnDemand</key>
			<dict>
				<key>ClientIdentifier</key>
				<string>com.apple.appstored</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.ondemandd.PeriodicLaunchServicesSyncActivity</key>
			<dict>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.appstored.CacheDelete</key>
		<true/>
		<key>com.apple.appstored.installcoordination</key>
		<true/>
		<key>com.apple.appstored.xpc</key>
		<true/>
		<key>com.apple.appstored.xpc.jobmanager</key>
		<true/>
		<key>com.apple.appstored.xpc.personalization</key>
		<true/>
		<key>com.apple.appstored.xpc.request</key>
		<true/>
		<key>com.apple.appstored.xpc.storequeue</key>
		<true/>
		<key>com.apple.appstored.xpc.updates</key>
		<true/>
		<key>com.apple.aps.appstored</key>
		<true/>
		<key>com.apple.ondemandd.CacheDelete</key>
		<true/>
		<key>com.apple.ondemandd.appinstall</key>
		<true/>
		<key>com.apple.ondemandd.appreview</key>
		<true/>
		<key>com.apple.ondemandd.client</key>
		<true/>
		<key>com.apple.ondemandd.devtools</key>
		<true/>
		<key>com.apple.ondemandd.launchservices</key>
		<true/>
		<key>com.apple.ondemandd.odrctl</key>
		<true/>
		<key>com.apple.private.alloy.appstore-idswake</key>
		<true/>
		<key>com.apple.storekitservice</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.AppStore</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AppStoreDaemon.framework/Support/appstored</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
