# diff/iOS/System/Library/LaunchDaemons/com.apple.filesystems.apfs_iosd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.filesystems.apfs_iosd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.system.powersources.timeremaining</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powersources.timeremaining</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.filesystems.apfs_iosd.ca_stats_report</key>
			<dict>
				<key>Interval</key>
				<integer>604800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresDeviceInactivity</key>
				<true/>
			</dict>
			<key>com.apple.filesystems.apfs_iosd.cryptex_stats_report</key>
			<dict>
				<key>CPUIntensive</key>
				<true/>
				<key>DiskIntensive</key>
				<true/>
				<key>Interval</key>
				<real>604800</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<false/>
				<key>RequiresDeviceInactivity</key>
				<true/>
			</dict>
			<key>com.apple.filesystems.apfs_iosd.fx_defrag</key>
			<dict>
				<key>CPUIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>259200</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.filesystems.apfs_iosd.hash_validation</key>
			<dict>
				<key>CPUIntensive</key>
				<true/>
				<key>DiskIntensive</key>
				<true/>
				<key>Interval</key>
				<real>604800</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresDeviceInactivity</key>
				<true/>
			</dict>
			<key>com.apple.filesystems.apfs_iosd.hint_extents</key>
			<dict>
				<key>CPUIntensive</key>
				<true/>
				<key>DiskIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>604800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresDeviceInactivity</key>
				<true/>
			</dict>
			<key>com.apple.filesystems.apfs_iosd.obj_cache_stats_report</key>
			<dict>
				<key>CPUIntensive</key>
				<true/>
				<key>DiskIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresDeviceInactivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Filesystems/apfs.fs/apfs_iosd</string>
	</array>
</dict>
</plist>
