# diff/iOS/System/Library/LaunchDaemons/com.apple.terminusd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>CFBundleIdentifier</key>
	<string>com.apple.terminusd</string>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>OS_ACTIVITY_ENABLE_DYNAMIC</key>
		<string>1</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/Library/Preferences/com.apple.terminusd.keepalive</key>
			<true/>
			<key>/var/mobile/Library/terminus/com.apple.terminusd.keepalive</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.terminusd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.networkrelay.launch.phs</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.networkrelay.launch.phs</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>com.apple.networkrelay.on-demand-setup</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.networkrelay.on-demand-setup</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.terminusd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/terminusd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
