# diff/iOS/System/Library/LaunchDaemons/com.apple.devicesharingd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>NSRunningFromLaunchd</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.devicesharingd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.accounts.changes</key>
		<dict>
			<key>AppleAccountAddedOrDeleted</key>
			<dict>
				<key>com.apple.account.AppleAccount</key>
				<dict>
					<key>ChangeTypes</key>
					<array>
						<integer>1</integer>
						<integer>3</integer>
					</array>
				</dict>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.kvs.store-did-change.com.apple.devicesharingd.deviceEnrollmentData</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.devicesharingd.deviceEnrollmentData</string>
			</dict>
			<key>com.vpg.managedassets.notify.iCloudSyncDown</key>
			<dict>
				<key>Notification</key>
				<string>com.vpg.managedassets.notify.iCloudSyncDown</string>
			</dict>
			<key>com.vpg.managedassets.notify.iCloudSyncUp</key>
			<dict>
				<key>Notification</key>
				<string>com.vpg.managedassets.notify.iCloudSyncUp</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>OnIncomingRequest</key>
			<dict>
				<key>serviceType</key>
				<string>DeviceSharingDaemonApplicationService</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.devicesharingd</key>
		<true/>
		<key>com.apple.devicesharingd.enrollmentassetservice</key>
		<true/>
		<key>com.apple.devicesharingd.guestuserremoteunlockservice</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/devicesharingd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
