# diff/iOS/System/Library/LaunchDaemons/com.apple.progressd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MallocMaxMagazines</key>
		<string>2</string>
		<key>MallocNanoMaxMagazines</key>
		<string>2</string>
	</dict>
	<key>Label</key>
	<string>com.apple.progressd</string>
	<key>LaunchEvents</key>
	<dict/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.progressd</key>
		<true/>
		<key>com.apple.progressd.aps</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/Frameworks/ClassKit.framework/progressd</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
