# diff/iOS/System/Library/LaunchDaemons/com.apple.nanoprefsyncd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AfterInitialDemand</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.nanoprefsyncd.2</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>NRPairedDeviceRegistryPairedDeviceDidChangeVersionDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.paireddevicedidchangeversion</string>
			</dict>
			<key>NRPairedDeviceRegistryWatchDidBecomeActiveDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.watchdidbecomeactive</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.nanoprefsync</key>
		<true/>
		<key>com.apple.nanoprefsyncd.first.pairedsync</key>
		<true/>
		<key>com.apple.nanoprefsyncd.pairedsync</key>
		<true/>
		<key>com.apple.private.alloy.preferencessync-idswake</key>
		<true/>
		<key>com.apple.private.alloy.preferencessync.pairedsync-idswake</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/NanoPreferencesSync.framework/nanoprefsyncd</string>
		<string>-companion</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
