# diff/iOS/System/Library/LaunchDaemons/com.apple.commandandcontrol.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.commandandcontrol</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.commandandcontrol</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.commandandcontrol.drag.xpc</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.speech.CommandAndControl.utility</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.commandandcontrol</key>
		<string></string>
	</dict>
	<key>NSHandsTrackingUsageDescription</key>
	<string>Hand tracking data includes hand skeleton, wrist, and forearm.</string>
	<key>NSWorldSensingUsageDescription</key>
	<string>World sensing data includes plane detection, scene geometry, and image tracking anchors.</string>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/CommandAndControl.app/CommandAndControl</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>ThrottleInterval</key>
	<integer>1</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
