# diff/iOS/System/Library/LaunchDaemons/com.apple.searchd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>NO_ASAN_DYLD_INSERT_LIBRARIES</key>
		<string>/usr/appleinternal/lib/sanitizers/libclang_rt.asan_ios_dynamic.dylib</string>
		<key>NO_ASAN_OPTIONS</key>
		<string>alloc_dealloc_mismatch=1:detect_stack_use_after_return=1:print_stats=1:atexit=1:help=1:verbosity=3:halt_on_error=1:strip_env=0</string>
		<key>NO_DYLD_FRAMEWORK_PATH</key>
		<string>/System/Library/PrivateFrameworks/MobileSpotlightIndex_asan.framework</string>
		<key>NO_DYLD_INSERT_LIBRARIES</key>
		<string>/usr/lib/libgmalloc.dylib</string>
		<key>NO_MallocScribble</key>
		<string>1</string>
		<key>NO_MallocStackLoggingNoCompacNo</key>
		<string>1</string>
		<key>NoMallocStackLogging</key>
		<string>1</string>
		<key>NoNSZombieEnabled</key>
		<string>1</string>
	</dict>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamMemoryLimit</key>
		<integer>4000</integer>
		<key>JetsamPriority</key>
		<integer>-49</integer>
	</dict>
	<key>Label</key>
	<string>com.apple.searchd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.searchd.duplicateOidCheck</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
			</dict>
			<key>com.apple.searchd.duplicateOidCheckAB</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>A</string>
			</dict>
			<key>com.apple.searchd.heartbeat</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>7200</integer>
				</dict>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>A</string>
			</dict>
			<key>com.apple.searchd.reportAppUsage</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>A</string>
			</dict>
			<key>com.apple.searchd.reportStorageUsage</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresNetworkConnectivity</key>
				<false/>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.corespotlight.developer.ReindexAllItems</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.corespotlight.developer.ReindexAllItems</string>
			</dict>
			<key>com.apple.corespotlight.developer.ReindexAllItemsWithIdentifiers</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.corespotlight.developer.ReindexAllItemsWithIdentifiers</string>
			</dict>
			<key>com.apple.datamigrator.migrationDidFinish</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.datamigrator.migrationDidFinish</string>
			</dict>
			<key>com.apple.spotlightui.prefschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.spotlightui.prefschanged</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.searchd.createEncryptedRankingBackup</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<real>86400</real>
				<key>GracePeriod</key>
				<integer>0</integer>
				<key>Interval</key>
				<integer>604800</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.searchd.defrag</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>2592000</integer>
				<key>GracePeriod</key>
				<integer>1296000</integer>
				<key>Interval</key>
				<integer>2592000</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.searchd.expirations</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>GracePeriod</key>
				<real>43200</real>
				<key>Interval</key>
				<real>86400</real>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.searchd.expirations.AB</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>GracePeriod</key>
				<real>43200</real>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresClassA</key>
				<true/>
			</dict>
			<key>com.apple.searchd.loadtrial</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>3600</integer>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>43200</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.searchd.sendmlfeedback</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>GracePeriod</key>
				<integer>0</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.searchd.updateRanking</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>GracePeriod</key>
				<real>43200</real>
				<key>Interval</key>
				<real>86400</real>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.MobileAsset.SpotlightResources</key>
		<true/>
		<key>com.apple.searchd</key>
		<true/>
		<key>com.apple.searchd.background</key>
		<true/>
		<key>com.apple.searchd.cachedelete</key>
		<true/>
		<key>com.apple.searchd.personaobserver</key>
		<true/>
		<key>com.apple.spotlight.IndexAgent</key>
		<true/>
		<key>com.apple.spotlight.IndexDelegateAgent</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.spotlight.SearchAgent</key>
		<true/>
	</dict>
	<key>NO_StandardErrorPath</key>
	<string>/var/mobile/stderr.log</string>
	<key>NO_StandardOutPath</key>
	<string>/var/mobile/stdout.log</string>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/Search.framework/searchd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>WaitForDebugger</key>
	<false/>
	<key>_ExponentialThrottling</key>
	<true/>
</dict>
</plist>
