# diff/iOS/System/Library/LaunchDaemons/com.apple.transparencyd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.transparencyd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.transparencyd.24h</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.transparencyd.dutycycle</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>28800</integer>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.transparencyd.self-validate</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>14400</integer>
				<key>NetworkTransferDirection</key>
				<string>Bidirectional</string>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.transparencyd</key>
		<true/>
		<key>com.apple.aps.transparencyd-coredata</key>
		<true/>
		<key>com.apple.transparencyd</key>
		<true/>
		<key>com.apple.transparencyd.accounts-support</key>
		<true/>
		<key>com.apple.transparencyd.ids-support</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/transparencyd</string>
</dict>
</plist>
