# diff/iOS/System/Library/LaunchDaemons/com.apple.ScreenTimeAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.ScreenTimeAgent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>_DKKnowledgeStorageDidTombstoneEventsNotification:AppUsage</key>
			<dict>
				<key>Name</key>
				<string>_DKKnowledgeStorageDidTombstoneEventsNotification</string>
				<key>Object</key>
				<string>/app/usage</string>
			</dict>
			<key>_DKKnowledgeStorageDidTombstoneEventsNotification:Notifications</key>
			<dict>
				<key>Name</key>
				<string>_DKKnowledgeStorageDidTombstoneEventsNotification</string>
				<key>Object</key>
				<string>/notification/usage</string>
			</dict>
			<key>_DKKnowledgeStorageDidTombstoneEventsNotification:WebUsage</key>
			<dict>
				<key>Name</key>
				<string>_DKKnowledgeStorageDidTombstoneEventsNotification</string>
				<key>Object</key>
				<string>/app/webUsage</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.dmd.iCloudAccount.didChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.dmd.iCloudAccount.didChange</string>
			</dict>
			<key>com.apple.family.family_updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.family.family_updated</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.mobile.lockdown.device_name_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.lockdown.device_name_changed</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.screentime.activity.analytics</key>
			<dict>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
			<key>com.apple.screentime.activity.app-collection</key>
			<dict>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.screentime.activity.cleanup</key>
			<dict>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.ScreenTimeAgent</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.Contacts</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.accessrequestor</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.accessresponder</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.ask</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.ask-for-time</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.command-line-tool</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.communication</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.diagnostics-service</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.downtime</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.exception</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.organization-status</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.persistence</key>
		<dict>
			<key>#Else</key>
			<true/>
			<key>#IfFeatureFlagEnabled</key>
			<string>ScreenTime/simplified_agent</string>
			<key>#Then</key>
			<false/>
		</dict>
		<key>com.apple.ScreenTimeAgent.private</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.reactor-tool</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.settings</key>
		<true/>
		<key>com.apple.ScreenTimeAgent.setup</key>
		<true/>
		<key>com.apple.aps.ScreenTime</key>
		<true/>
		<key>com.apple.legacyagent.usersettingsprovider</key>
		<true/>
		<key>com.apple.legacyagent.usersettingsstore</key>
		<true/>
		<key>com.apple.private.alloy.digitalhealth-idswake</key>
		<true/>
		<key>com.apple.private.alloy.screentime-idswake</key>
		<true/>
		<key>com.apple.private.alloy.screentimelocal-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.ScreenTimeDowntimeNotifications</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.ScreenTimeNotifications</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/ScreenTimeCore.framework/ScreenTimeAgent</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
