# diff/iOS/System/Library/LaunchDaemons/com.apple.fileindexerd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<dict>
		<key>#IfFeatureFlagDisabled</key>
		<string>FileIndexer/IndexerDaemon</string>
		<key>#Then</key>
		<true/>
	</dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LIBDISPATCH_STRICT</key>
		<string>YES</string>
		<key>OBJC_DEBUG_MISSING_POOLS</key>
		<string>halt</string>
	</dict>
	<key>Label</key>
	<string>com.apple.fileindexerd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.fsevents.matching</key>
		<dict>
			<key>com.apple.fileindexerd_application.fsevents</key>
			<dict>
				<key>Path</key>
				<string>/private/var/mobile/Containers/Shared/AppGroup</string>
			</dict>
			<key>com.apple.fileindexerd_local.fsevents</key>
			<dict>
				<key>Path</key>
				<string>/private/var/mobile/Containers/Data/Application</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.FileIndexer</key>
		<true/>
		<key>com.apple.corespotlight.daemon.fileindexer</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/FileIndexerDaemon.framework/Support/fileindexerd</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
