# diff/iOS/System/Library/LaunchDaemons/com.apple.familycircled.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>CFBundleIdentifier</key>
	<string>com.apple.familycircled</string>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentalVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.familycircled</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.family.cache.refresh.backgroundTask</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>84600</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.family.heartbeat.backgroundTask</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>82800</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.family.osupdate.backgroundTask</key>
			<dict>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>__ABDataBaseChangedByOtherProcessNotification</key>
			<dict>
				<key>Notification</key>
				<string>__ABDataBaseChangedByOtherProcessNotification</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.family</key>
		<true/>
		<key>com.apple.aps.familycircled.cloudkit.familycircle.private</key>
		<true/>
		<key>com.apple.aps.familycircled.cloudkit.familycircle.shared</key>
		<true/>
		<key>com.apple.aps.familycircled.cloudkit.familycircleprivate.private</key>
		<true/>
		<key>com.apple.aps.familycircled.cloudkit.familycircleprivate.shared</key>
		<true/>
		<key>com.apple.family.ageRange.xpc</key>
		<true/>
		<key>com.apple.family.notifier</key>
		<true/>
		<key>com.apple.familycircle.agent</key>
		<true/>
		<key>com.apple.familycircled.settings</key>
		<true/>
		<key>com.apple.familycircled.sharing</key>
		<true/>
		<key>com.apple.familycircled.testing</key>
		<true/>
		<key>com.apple.private.alloy.familycircle-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/FamilyCircle.framework/familycircled</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>3</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
