# diff/iOS/System/Library/LaunchDaemons/com.apple.SensorKitALSHelper.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.SensorKitALSHelper</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>ALS recording status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SensorKit.als</string>
			</dict>
			<key>Acoustic Settings recording status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SensorKit.hearing.acousticSettings</string>
			</dict>
			<key>Media events recording status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SensorKit.mediaEvents</string>
			</dict>
			<key>Messages Usage Report recording status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SensorKit.messagesUsageReport</string>
			</dict>
			<key>Phone Usage Report recording status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SensorKit.phoneUsageReport</string>
			</dict>
			<key>Sleep Sessions recording status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SensorKit.sleep.sessions</string>
			</dict>
			<key>Stride calibration recording status </key>
			<dict>
				<key>Notification</key>
				<string>com.apple.private.SensorKit.pedometer.stridecalibration</string>
			</dict>
			<key>Usage Report recording status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SensorKit.deviceUsageReport</string>
			</dict>
			<key>Visits recording status </key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SensorKit.visits</string>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/SensorKitALSHelper</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
