# diff/iOS/System/Library/LaunchDaemons/com.apple.matd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.matd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.purplebuddy.setupexited</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.purplebuddy.setupexited</string>
			</dict>
			<key>com.apple.welcomekitinternalsettings.dismissed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.welcomekitinternalsettings.dismissed</string>
			</dict>
		</dict>
		<key>com.apple.nsurlsessiond-events</key>
		<dict>
			<key>LaunchOnDemand</key>
			<dict>
				<key>ClientIdentifier</key>
				<string>com.apple.matd</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.matd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/WelcomeKit.framework/matd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
