# diff/iOS/System/Library/LaunchDaemons/com.apple.akd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.akd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.LocalAuthentication.ratchet.StateDidChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LocalAuthentication.ratchet.StateDidChange</string>
			</dict>
			<key>com.apple.managedconfiguration.passcodechanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.passcodechanged</string>
			</dict>
			<key>com.apple.mobile.keybagd.lock_status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.lock_status</string>
			</dict>
			<key>com.apple.security.kcsharing.groupsupdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.security.kcsharing.groupsupdated</string>
			</dict>
			<key>com.apple.springboard.lockstate</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.lockstate</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>AppSignInEndorser</key>
			<dict>
				<key>serviceType</key>
				<string>RPUserNotificationAppSignIn</string>
				<key>type</key>
				<string>discovery</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.ak.hb</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.ak.listupdate</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<real>3600</real>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.ak.passkey.healthcheck</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.ak.postinstall</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<real>3600</real>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.ak.signinwithapple.healthcheck</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>180</integer>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresClassA</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.ak.anisette.xpc</key>
		<true/>
		<key>com.apple.ak.appleidpasskey.xpc</key>
		<true/>
		<key>com.apple.ak.aps</key>
		<true/>
		<key>com.apple.ak.auth.xpc</key>
		<true/>
		<key>com.apple.ak.authorizationservices.xpc</key>
		<true/>
		<key>com.apple.ak.custodian.xpc</key>
		<true/>
		<key>com.apple.ak.inheritance.xpc</key>
		<true/>
		<key>com.apple.ak.privateemail.xpc</key>
		<true/>
		<key>com.apple.ak.puffin.xpc</key>
		<true/>
		<key>com.apple.ak.signinwithapple.xpc</key>
		<true/>
		<key>com.apple.ak.symmetrickey.xpc</key>
		<true/>
		<key>com.apple.ak.walrus.xpc</key>
		<true/>
		<key>com.apple.private.alloy.anisette-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AuthKit.framework/akd</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>3</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
