# diff/iOS/System/Library/LaunchDaemons/com.apple.corespeechd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/var/mobile/Library/VoiceTrigger/KeepAlive</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.corespeechd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.siri.bg_system_task.daily-asset-subscription-cleanup</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.siri.bg_system_task.daily-speech-profile-maintenance</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
			</dict>
			<key>com.apple.siri.bg_system_task.post-install-asset-subscription</key>
			<dict>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
			</dict>
			<key>com.apple.siri.bg_system_task.post-install-speech-profile-migration</key>
			<dict>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>kBGSTProtectionClassC</string>
				<key>TrySchedulingBefore</key>
				<integer>300</integer>
			</dict>
			<key>com.apple.siri.bg_system_task.replay-record-pruning</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>900</integer>
				</dict>
			</dict>
		</dict>
		<key>com.apple.cascade.xpc.setChange</key>
		<dict>
			<key>com.apple.corespeechd</key>
			<dict/>
		</dict>
		<key>com.apple.cascade.xpc_event.setChange</key>
		<dict>
			<key>com.apple.corespeechd</key>
			<dict/>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>Application Uninstalled</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.mobile.keybagd.first_unlock</key>
		<dict>
			<key>Name</key>
			<string>com.apple.mobile.keybagd.first_unlock</string>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AFLanguageCodeDidChangeDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>AFLanguageCodeDidChangeDarwinNotification</string>
			</dict>
			<key>VT Phrase Type changed</key>
			<dict>
				<key>Notification</key>
				<string>VT Phrase Type changed</string>
			</dict>
			<key>com.apple.MobileAsset.AutoAssetNotification^com.apple.MobileAsset.MAAutoAsset^STARTUP_ACTIVATED</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.AutoAssetNotification^com.apple.MobileAsset.MAAutoAsset^STARTUP_ACTIVATED</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssets.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssets.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssets.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssets.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsIPad.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsIPad.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsIPad.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsIPad.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsMarsh.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsMarsh.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsMarsh.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsMarsh.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsStudioDisplay.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsStudioDisplay.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsStudioDisplay.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsStudioDisplay.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsWatch.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsWatch.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsWatch.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsWatch.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerHSAssets.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerHSAssets.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerHSAssets.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerHSAssets.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerHSAssetsIPad.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerHSAssetsIPad.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerHSAssetsIPad.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerHSAssetsIPad.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerHSAssetsWatch.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerHSAssetsWatch.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerHSAssetsWatch.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerHSAssetsWatch.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.assistant.domain.didChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.domain.didChange</string>
			</dict>
			<key>com.apple.assistant.domain.preferences.didChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.domain.preferences.didChange</string>
			</dict>
			<key>com.apple.assistant.domain.priority.didChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.domain.priority.didChange</string>
			</dict>
			<key>com.apple.audio.AOP.enable</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.audio.AOP.enable</string>
			</dict>
			<key>com.apple.bluetooth.WirelessSplitterOn</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.bluetooth.WirelessSplitterOn</string>
			</dict>
			<key>com.apple.coreaudio.BorealisToggled</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.coreaudio.BorealisToggled</string>
			</dict>
			<key>com.apple.coreaudio.speechDetectionVAD.created</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.coreaudio.speechDetectionVAD.created</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.siri.understanding</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.siri.understanding</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.speech.automaticspeechrecognition</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.speech.automaticspeechrecognition</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.voiceassistant.perception</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.voiceassistant.perception</string>
			</dict>
			<key>com.apple.springboard.finishedstartup</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.finishedstartup</string>
			</dict>
			<key>com.apple.sysdiagnose.sysdiagnoseStarted</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sysdiagnose.sysdiagnoseStarted</string>
			</dict>
			<key>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_ASR_ASSISTANT</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_ASR_ASSISTANT</string>
			</dict>
			<key>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_ATTENTION_ASSETS</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_ATTENTION_ASSETS</string>
			</dict>
			<key>com.apple.voicetrigger.RemoteDarwin.ConnectionChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicetrigger.RemoteDarwin.ConnectionChanged</string>
			</dict>
			<key>com.apple.voicetrigger.RemoteDarwin.EarlyDetect</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicetrigger.RemoteDarwin.EarlyDetect</string>
			</dict>
			<key>com.apple.voicetrigger.XPCRestarted</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicetrigger.XPCRestarted</string>
			</dict>
			<key>com.apple.voicetrigger.enablePolicyChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicetrigger.enablePolicyChanged</string>
			</dict>
			<key>kAFPreferencesDidChangeDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>kAFPreferencesDidChangeDarwinNotification</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.assistant.domain.audio.speech.service</key>
		<true/>
		<key>com.apple.assistant.domain.debug.speech.service</key>
		<true/>
		<key>com.apple.assistant.domain.keyworddetector.speech.service</key>
		<true/>
		<key>com.apple.corespeech.audioinjection.plugin.xpc</key>
		<true/>
		<key>com.apple.corespeech.audioinjection.xpc</key>
		<true/>
		<key>com.apple.corespeech.benchmark.xpc</key>
		<true/>
		<key>com.apple.corespeech.corespeechd.activation.xpc</key>
		<true/>
		<key>com.apple.corespeech.corespeechd.attending.service</key>
		<true/>
		<key>com.apple.corespeech.corespeechd.attsiri.service</key>
		<true/>
		<key>com.apple.corespeech.corespeechd.endpointer.service</key>
		<true/>
		<key>com.apple.corespeech.corespeechd.rchandling.service</key>
		<true/>
		<key>com.apple.corespeech.corespeechd.speechprofileadmin.service</key>
		<true/>
		<key>com.apple.corespeech.corespeechd.uaapservice</key>
		<true/>
		<key>com.apple.corespeech.corespeechd.voiceid.xpc</key>
		<true/>
		<key>com.apple.corespeech.corespeechd.xpc</key>
		<true/>
		<key>com.apple.corespeech.corespeechservices</key>
		<true/>
		<key>com.apple.corespeech.voicetriggerservice</key>
		<true/>
		<key>com.apple.siri.attendingstates.xpc</key>
		<true/>
		<key>com.apple.siri.audio_message_service.xpc</key>
		<true/>
		<key>com.apple.siri.audiopowerupdate.xpc</key>
		<true/>
		<key>com.apple.siri.localspeechrecognitionbridge.xpc</key>
		<true/>
		<key>com.apple.siri.ssrvtuitrainingservice.xpc</key>
		<true/>
		<key>com.apple.ssr.rpisamplingservice</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CoreSpeech.framework/corespeechd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>_Conclave</key>
	<string>com.apple.corespeechd.conclave</string>
</dict>
</plist>
