# diff/iOS/System/Library/LaunchDaemons/com.apple.voiced.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.voiced</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.MobileAsset.VoiceServices.GryphonVoice.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceServices.GryphonVoice.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceServices.VoiceResources.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceServices.VoiceResources.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceServicesVocalizerVoice.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceServicesVocalizerVoice.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.voiceservices.trigger.asset-force-update</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voiceservices.trigger.asset-force-update</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.voiced</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.voiced.CacheDelete</key>
		<true/>
		<key>com.apple.voiceservices.keepalive</key>
		<true/>
		<key>com.apple.voiceservices.tts</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/VoiceServices.framework/Support/voiced</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
