# beta/macOS/System/Library/Templates/Data/Library/Apple/System/Library/LaunchDaemons/com.apple.XProtect.daemon.scan.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CHECK_IN</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.XProtect.daemon.scan</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.XProtect.PluginService.daemon.fast.scan</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>CPUIntensive</key>
				<true/>
				<key>DiskIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>21600</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.XProtect.PluginService.daemon.scan</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>CPUIntensive</key>
				<true/>
				<key>DiskIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.XProtect.PluginService.daemon.slow.scan</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>CPUIntensive</key>
				<true/>
				<key>DiskIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>604800</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Background</string>
	<key>Program</key>
	<string>/Library/Apple/System/Library/CoreServices/XProtect.app/Contents/MacOS/XProtect</string>
</dict>
</plist>
