# beta/System/Library/NanoLaunchDaemons/com.apple.screenshotsyncd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>Label</key>
	<string>com.apple.screenshotsyncd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.private.alloy.screenshotter-idswake</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/screenshotsyncd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
	<key>WorkingDirectory</key>
	<string>/var/mobile/Library/Caches</string>
</dict>
</plist>
