# beta/System/Library/NanoLaunchDaemons/com.apple.appconduitd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/private/var/mobile/Library/AppConduit/install_in_progress</key>
			<true/>
			<key>/private/var/mobile/Library/AppConduit/local_appstore_install_in_progress</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.appconduitd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>MISUPPTrustRevoked</key>
			<dict>
				<key>Notification</key>
				<string>MISUPPTrustRevoked</string>
			</dict>
			<key>MISUPPTrustSet</key>
			<dict>
				<key>Notification</key>
				<string>MISUPPTrustSet</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.pairedsync.syncDidComplete</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.pairedsync.syncDidComplete</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.appconduitd.appstate.pairedsync</key>
		<true/>
		<key>com.apple.appconduitd.device-connection</key>
		<true/>
		<key>com.apple.appconduitd.pairedsync</key>
		<true/>
		<key>com.apple.private.alloy.appconduit-idswake</key>
		<true/>
		<key>com.apple.private.alloy.appconduit.v2-idswake</key>
		<true/>
		<key>com.apple.private.alloy.appsyncconduit-idswake</key>
		<true/>
		<key>com.apple.private.alloy.appsyncconduit.v2-idswake</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/AppConduit.framework/Support/appconduitd</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
