# beta/System/Library/NanoLaunchDaemons/com.apple.companionappd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.companionappd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.pairedsync.syncDidComplete</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.pairedsync.syncDidComplete</string>
			</dict>
			<key>com.apple.sockpuppet.activeComplicationsPreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sockpuppet.activeComplicationsPreferencesChangedNotification</string>
			</dict>
			<key>com.apple.sockpuppet.idslaunchnotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sockpuppet.idslaunchnotification</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.companionappd</key>
		<true/>
		<key>com.apple.companionappd.pingsubscriber</key>
		<true/>
		<key>com.apple.private.alloy.sockpuppet-idswake</key>
		<true/>
		<key>com.apple.private.alloy.sockpuppet.classd-idswake</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/Frameworks/WatchKit.framework/Support/companionappd</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
