# beta/System/Library/NanoLaunchDaemons/com.apple.nanoappregistryd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.nanoappregistryd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>ApplicationsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LaunchServices.ApplicationsChanged</string>
			</dict>
			<key>IDS</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appregistrysync.idslaunchnotification</string>
			</dict>
			<key>LockState</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.lockstate</string>
			</dict>
			<key>Pairing</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.nanoappregistry.pairedsync</key>
		<true/>
		<key>com.apple.nanoappregistry.workspace</key>
		<true/>
		<key>com.apple.private.alloy.appregistrysync-idswake</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/NanoAppRegistry.framework/Support/nanoappregistryd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
