# beta/System/Library/NanoLaunchDaemons/com.apple.eventkitsyncd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.eventkitsyncd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CalDefaultCalendarChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>CalDefaultCalendarChangedNotification</string>
			</dict>
			<key>CalEventDeclineEventPrefChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilecal.showdeclinedchanged</string>
			</dict>
			<key>CalEventShowInviteeDeclinesPrefChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilecal.invitationalertschanged</string>
			</dict>
			<key>EventKitSyncTesting</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.eventkitsync.fullsync</string>
			</dict>
			<key>IDS</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.eventkitsync.idslaunchnotification</string>
			</dict>
			<key>IDS2</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.eventkitmutation.idslaunchnotification</string>
			</dict>
			<key>LastDeselectedCalendarSyncHashes</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilecal.preference.notification.selectedCalendarSyncHashesChanged</string>
			</dict>
			<key>SignificantTimeChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>SignificantTimeChangeNotification</string>
			</dict>
			<key>_CalDatabaseChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>_CalDatabaseChangedNotification</string>
			</dict>
			<key>com.apple.calendar.database.preference.notification.enableTravelAdvisoriesForAutomaticBehavior</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.calendar.database.preference.notification.enableTravelAdvisoriesForAutomaticBehavior</string>
			</dict>
			<key>com.apple.calendar.defaultAlarmChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.calendar.defaultAlarmChangedNotification</string>
			</dict>
			<key>com.apple.eventkitsync.analytics.flush</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.eventkitsync.analytics.flush</string>
			</dict>
			<key>com.apple.eventkitsync.force</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.eventkitsync.force</string>
			</dict>
			<key>com.apple.eventkitsync.zapversion</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.eventkitsync.zapversion</string>
			</dict>
			<key>com.apple.mobilecal.preference.notification.showweeknumbers</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilecal.preference.notification.showweeknumbers</string>
			</dict>
			<key>com.apple.reminderkit.storeChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.reminderkit.storeChanged</string>
			</dict>
			<key>kCalTimeZonePrefChangedDistributedNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilecal.timezonechanged</string>
			</dict>
			<key>overlayCalendarID</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilecal.preference.notification.overlayCalendarID</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.eventkitsync.analytics.daily</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.EventKitSyncServices.server</key>
		<true/>
		<key>com.apple.eventkitsync.pairedsync</key>
		<true/>
		<key>com.apple.internal.eventkitsync.testcontrol</key>
		<true/>
		<key>com.apple.private.alloy.eventkitmutation-idswake</key>
		<true/>
		<key>com.apple.private.alloy.eventkitsync-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.eventkitsync.notifications</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/eventkitsyncd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
