# beta/System/Library/NanoLaunchDaemons/com.apple.pairedsyncd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/var/mobile/Library/PairedSyncKeepAlive</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.pairedsyncd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>MobileKeybagFirstUnlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>NanoRegistryDevicePaired</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
			<key>NanoRegistrySyncAdvice</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.pairedsync.isokaytosync</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.pairedsyncd.connection</key>
		<true/>
		<key>com.apple.pairedsyncd.progress</key>
		<true/>
		<key>com.apple.pairedsyncd.syncstate</key>
		<true/>
		<key>com.apple.pairedsyncd.tool</key>
		<true/>
	</dict>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/PairedSync.framework/Support/pairedsyncd</string>
	<key>RunAtLoad</key>
	<true/>
	<key>TransactionTimeLimitEnabled</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
