# beta/iOS/System/Library/LaunchDaemons/com.apple.nand_task_scheduler.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>BGTaskSchedulerPermittedIdentifiers</key>
	<array>
		<string>com.apple.idleStack_low</string>
		<string>com.apple.idleStack_lowNoDI</string>
		<string>com.apple.idleStack_lowNoSUI</string>
		<string>com.apple.idleStack_med</string>
		<string>com.apple.idleStack_high</string>
	</array>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.nand_task_scheduler</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.idleStack_high</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>3600</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>2400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresSignificantUserInactivity</key>
				<false/>
				<key>ResourceIntensive</key>
				<true/>
				<key>Resources</key>
				<integer>4</integer>
			</dict>
			<key>com.apple.idleStack_low</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>69120</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresSignificantUserInactivity</key>
				<true/>
				<key>ResourceIntensive</key>
				<true/>
				<key>Resources</key>
				<integer>4</integer>
			</dict>
			<key>com.apple.idleStack_lowNoDI</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86460</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>69120</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresSignificantUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.idleStack_lowNoSUI</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86460</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>69120</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>ResourceIntensive</key>
				<true/>
				<key>Resources</key>
				<integer>4</integer>
			</dict>
			<key>com.apple.idleStack_med</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>21600</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>16200</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresSignificantUserInactivity</key>
				<false/>
				<key>ResourceIntensive</key>
				<true/>
				<key>Resources</key>
				<integer>4</integer>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.nand_task_scheduler</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/NANDTaskScheduler</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/NANDTaskScheduler</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
