# beta/iOS/System/Library/LaunchDaemons/com.apple.homed.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AppIntents</key>
	<dict>
		<key>MetadataAbsolutePaths</key>
		<array>
			<string>/System/Library/PrivateFrameworks/HomeKitDaemon.framework/</string>
			<string>/System/Library/PrivateFrameworks/HomeKitDaemonLegacy.framework</string>
		</array>
	</dict>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/var/mobile/Library/homed/com.apple.relaunch.assertion.plist</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.homed</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.accounts.changes</key>
		<dict>
			<key>AppleAccountUpdated</key>
			<dict>
				<key>com.apple.account.AppleAccount</key>
				<dict>
					<key>ChangeTypes</key>
					<array>
						<integer>1</integer>
						<integer>2</integer>
						<integer>3</integer>
					</array>
				</dict>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.CascadeSets.DonateNow</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.CascadeSets.DonateNow</string>
			</dict>
		</dict>
		<key>com.apple.passd.matching</key>
		<dict>
			<key>PKPassLibraryCanAddPassChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>PKPassLibraryCanAddPassChangeNotification</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>com.apple.home.messaging</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.home.messaging</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.homekit</key>
		<true/>
		<key>com.apple.corespotlight.daemon.homekit</key>
		<true/>
		<key>com.apple.homed.aps</key>
		<true/>
		<key>com.apple.homed.xpc</key>
		<true/>
		<key>com.apple.homed.xpc.accessory-setup</key>
		<true/>
		<key>com.apple.homekit.coredata.xpc</key>
		<true/>
		<key>com.apple.homekit.pairedsync</key>
		<true/>
		<key>com.apple.homekit.xpc</key>
		<true/>
		<key>com.apple.matter.framework.xpc</key>
		<true/>
		<key>com.apple.matter.native.xpc</key>
		<true/>
		<key>com.apple.matter.support.xpc</key>
		<true/>
		<key>com.apple.nfcd.xpc.homed.uaevents</key>
		<true/>
		<key>com.apple.private.alloy.home-idswake</key>
		<true/>
		<key>com.apple.private.alloy.home.invite-idswake</key>
		<true/>
		<key>com.apple.private.alloy.willow-idswake</key>
		<true/>
		<key>com.apple.private.alloy.willow.proxy-idswake</key>
		<true/>
		<key>com.apple.private.alloy.willow.stream-idswake</key>
		<true/>
		<key>com.apple.private.appintents.delegate.com.apple.homed</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.Home</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/HomeKitDaemon.framework/Support/homed</string>
	</array>
	<key>PublishesEvents</key>
	<dict>
		<key>com.apple.HMModernMessaging.launch</key>
		<dict>
			<key>RequireTrusted</key>
			<true/>
		</dict>
	</dict>
	<key>UserName</key>
	<string>mobile</string>
	<key>_ExponentialThrottling</key>
	<true/>
</dict>
</plist>
