# beta/iOS/System/Library/LaunchDaemons/com.apple.assessmentagent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/var/mobile/Library/com.apple.assessmentagent/PolicyStore</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.assessmentagent</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.assessmentagent</key>
		<true/>
		<key>com.apple.assessmentagent.activeRestrictionUUIDFetching</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/usr/libexec/assessmentagent</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
