# beta/iOS/System/Library/LaunchDaemons/com.apple.fontservicesd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.fontservicesd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>FontServicesDaemonCheckSubscriptionNotification</key>
			<dict>
				<key>Name</key>
				<string>FontServicesDaemonCheckSubscriptionNotification</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.MobileAsset.CoreTextAssets.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.CoreTextAssets.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.CoreTextAssets.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.CoreTextAssets.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.Font8.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.Font8.ma.cached-metadata-updated</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.fontservicesd.check-asset-update</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>ExpectedDuration</key>
				<integer>15</integer>
				<key>GracePeriod</key>
				<integer>43200</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.fontservicesd.subscription-support</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>ExpectedDuration</key>
				<integer>10</integer>
				<key>GracePeriod</key>
				<integer>21600</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.fontservicesd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/FontServices.framework/Support/fontservicesd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
