# beta/iOS/System/Library/LaunchDaemons/com.apple.searchtoold.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>NO_MallocScribble</key>
		<string>1</string>
		<key>NO_MallocStackLoggingNoCompacNo</key>
		<string>1</string>
		<key>NoMallocStackLogging</key>
		<string>1</string>
		<key>NoNSZombieEnabled</key>
		<string>1</string>
	</dict>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamMemoryLimit</key>
		<integer>4000</integer>
		<key>JetsamPriority</key>
		<integer>-49</integer>
	</dict>
	<key>Label</key>
	<string>com.apple.searchtoold</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.omniSearch.answerSynthesis</key>
		<true/>
		<key>com.apple.omniSearch.echo</key>
		<true/>
		<key>com.apple.omniSearch.personalSearch</key>
		<true/>
		<key>com.apple.omniSearch.search</key>
		<true/>
	</dict>
	<key>NO_StandardErrorPath</key>
	<string>/var/mobile/stderr.log</string>
	<key>NO_StandardOutPath</key>
	<string>/var/mobile/stdout.log</string>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/OmniSearch.framework/searchtoold</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
