# beta/iOS/System/Library/LaunchDaemons/com.apple.icloud.searchpartyd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.icloud.searchpartyd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bluetooth.discovery</key>
		<dict>
			<key>ProximityPairing</key>
			<dict>
				<key>bleRSSIThresholdHint</key>
				<integer>-50</integer>
				<key>discoveryFlags</key>
				<array>
					<string>ObjectSetup</string>
					<string>PassiveRSSI</string>
				</array>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>FMAccessoryOnboardingDidAppear</key>
			<dict>
				<key>Name</key>
				<string>FMAccessoryOnboardingDidAppear</string>
			</dict>
			<key>FMAccessoryOnboardingDidDisappear</key>
			<dict>
				<key>Name</key>
				<string>FMAccessoryOnboardingDidDisappear</string>
			</dict>
			<key>PNPParingVCDidAppearNotification</key>
			<dict>
				<key>Name</key>
				<string>PNPParingVCDidAppearNotification</string>
			</dict>
			<key>PNPParingVCDidDisappearNotification</key>
			<dict>
				<key>Name</key>
				<string>PNPParingVCDidDisappearNotification</string>
			</dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CKIdentityUpdateNotification</key>
			<dict>
				<key>Notification</key>
				<string>CKIdentityUpdateNotification</string>
			</dict>
			<key>FMLDevicesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>FMLDevicesChangedNotification</string>
			</dict>
			<key>FMLFollowersChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>FMLFollowersChangedNotification</string>
			</dict>
			<key>FMLMeDeviceChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>FMLMeDeviceChangedNotification</string>
			</dict>
			<key>RTLocationsOfInterestDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>RTLocationsOfInterestDidChangeNotification</string>
			</dict>
			<key>com.apple.AOSNotification.FMIPStateDidChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AOSNotification.FMIPStateDidChange</string>
			</dict>
			<key>com.apple.FindMyDevice.LocatableStateReported</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.FindMyDevice.LocatableStateReported</string>
			</dict>
			<key>com.apple.bluetooth.pairingWithReason</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.bluetooth.pairingWithReason</string>
			</dict>
			<key>com.apple.family.family_updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.family.family_updated</string>
			</dict>
			<key>com.apple.icloud.findmydeviced.localActivationLockInfoChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.findmydeviced.localActivationLockInfoChanged</string>
			</dict>
			<key>com.apple.locationd.authorization</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.authorization</string>
			</dict>
			<key>com.apple.locationd/Prefs</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd/Prefs</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.icloud.searchpartyuseragent.post-install</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.searchpartyd</key>
		<true/>
		<key>com.apple.findmy.push.com.apple.findmy.localfindable</key>
		<true/>
		<key>com.apple.icloud.searchparty.delegatedlocationpublish.items</key>
		<true/>
		<key>com.apple.icloud.searchparty.locationfetch.items</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.accessorydiscoverysession</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.advertisementcache</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.aps</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.beaconmanager</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.beaconmanager.simplebeacon</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.beaconsharingservice</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.finderstatemanager</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.intentsession</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.localFindableConnectionMaterialSession</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.ownersession</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.pairingmanager</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.scheduler</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.securelocations</key>
		<true/>
		<key>com.apple.icloud.searchpartyd.unknowndiscoverysession</key>
		<true/>
		<key>com.apple.private.alloy.findmy.itemsharing-crossaccount-idswake</key>
		<true/>
		<key>com.apple.private.alloy.fmd.local-idswake</key>
		<true/>
		<key>com.apple.searchparty.BeaconScanningSession</key>
		<true/>
		<key>com.apple.searchparty.btfindingsession</key>
		<true/>
		<key>com.apple.searchparty.managedperipheral</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.FindMySafetyAlertsNotifications</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/searchpartyd</string>
	<key>Umask</key>
	<integer>23</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
