# beta/iOS/System/Library/LaunchDaemons/com.apple.calaccessd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.calaccessd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CalSyncClientBeginningMultiSave</key>
			<dict>
				<key>Notification</key>
				<string>CalSyncClientBeginningMultiSave</string>
			</dict>
			<key>CalSyncClientFinishedMultiSave</key>
			<dict>
				<key>Notification</key>
				<string>CalSyncClientFinishedMultiSave</string>
			</dict>
			<key>EKNotificationCountChangedExternallyNotification</key>
			<dict>
				<key>Notification</key>
				<string>EKNotificationCountChangedExternallyNotification</string>
			</dict>
			<key>SignificantTimeChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>SignificantTimeChangeNotification</string>
			</dict>
			<key>_CalDatabaseChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>_CalDatabaseChangedNotification</string>
			</dict>
			<key>_CalDatabaseIntegrationDataChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>_CalDatabaseIntegrationDataChangedNotification</string>
			</dict>
			<key>com.apple.accessibility.prefers.horizontal.text</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.prefers.horizontal.text</string>
			</dict>
			<key>com.apple.appprotection.change.hidden</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appprotection.change.hidden</string>
			</dict>
			<key>com.apple.appprotection.change.locked</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appprotection.change.locked</string>
			</dict>
			<key>com.apple.calendar.database.preference.notification.suggestEventLocations</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.calendar.database.preference.notification.suggestEventLocations</string>
			</dict>
			<key>com.apple.eventkit.preference.notification.AlertInviteeDeclines</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.eventkit.preference.notification.AlertInviteeDeclines</string>
			</dict>
			<key>com.apple.eventkit.preference.notification.UnselectedCalendarIdentifiersForFocusMode</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.eventkit.preference.notification.UnselectedCalendarIdentifiersForFocusMode</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.mobilecal.invitationalertschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilecal.invitationalertschanged</string>
			</dict>
			<key>com.apple.mobilecal.preference.notification.calendarsExcludedFromNotifications</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilecal.preference.notification.calendarsExcludedFromNotifications</string>
			</dict>
			<key>com.apple.mobilecal.preference.notification.overlayCalendarID</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilecal.preference.notification.overlayCalendarID</string>
			</dict>
			<key>com.apple.mobilecal.preference.notification.weekStart</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilecal.preference.notification.weekStart</string>
			</dict>
			<key>com.apple.mobilecal.timezonechanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilecal.timezonechanged</string>
			</dict>
			<key>com.apple.navd.wakeUpForHypothesisUpdate</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.navd.wakeUpForHypothesisUpdate</string>
			</dict>
			<key>com.apple.reminderkit.storeChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.reminderkit.storeChanged</string>
			</dict>
			<key>com.apple.springboard.finishedstartup</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.finishedstartup</string>
			</dict>
			<key>kCalBirthdayDefaultAlarmChangedNote</key>
			<dict>
				<key>Notification</key>
				<string>kCalBirthdayDefaultAlarmChangedNote</string>
			</dict>
			<key>kCalEventOccurrenceCacheChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>kCalEventOccurrenceCacheChangedNotification</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.calendar.daemon.analytics</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<real>86400</real>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.calendar.daemon.attachmentcleanup</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<real>86400</real>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.calendar.daemon.changetablecleanup</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<real>28800</real>
				<key>GracePeriod</key>
				<real>28800</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.calendar.daemon.conferencerenewal</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>259200</integer>
				<key>GracePeriod</key>
				<integer>259200</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.calendar.daemon.databasecleanup</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<real>86400</real>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.calendar.daemon.occurrencecacheupdate</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<real>86400</real>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.CalendarExtensionSupport</key>
		<true/>
		<key>com.apple.calaccessd</key>
		<true/>
		<key>com.apple.calaccessd.xpc</key>
		<dict>
			<key>HideUntilCheckIn</key>
			<true/>
		</dict>
		<key>com.apple.corespotlight.daemon.calendar</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.mobilecal</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CalendarDaemon.framework/Support/calaccessd</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>4</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
