# beta/iOS/System/Library/LaunchDaemons/com.apple.mobile.lockdown.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.mobile.lockdown</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.lockdown.KIS</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>Active</key>
					<true/>
					<key>TransportType</key>
					<integer>11</integer>
				</dict>
				<key>IOProviderClass</key>
				<string>IOPortTransportState</string>
			</dict>
			<key>com.apple.lockdown.USB2.Device</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>Active</key>
					<true/>
					<key>DataRole</key>
					<integer>1</integer>
					<key>TransportType</key>
					<integer>2</integer>
				</dict>
				<key>IOProviderClass</key>
				<string>IOPortTransportState</string>
			</dict>
			<key>com.apple.lockdown.USB3.Device</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>Active</key>
					<true/>
					<key>DataRole</key>
					<integer>1</integer>
					<key>TransportType</key>
					<integer>3</integer>
				</dict>
				<key>IOProviderClass</key>
				<string>IOPortTransportState</string>
			</dict>
			<key>com.apple.lockdown.USBHost</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>AppleUSBCableType</key>
					<string>USBHost</string>
				</dict>
				<key>IOProviderClass</key>
				<string>IOResources</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
			<key>com.apple.datamigrator.migrationDidFinish</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.datamigrator.migrationDidFinish</string>
			</dict>
			<key>com.apple.mobile.disk_image_mounted</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.disk_image_mounted</string>
			</dict>
			<key>com.apple.mobile.storage_unmounted</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.storage_unmounted</string>
			</dict>
		</dict>
	</dict>
	<key>MinimalBootProfile</key>
	<true/>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/lockdownd</string>
	</array>
	<key>RemoteServices</key>
	<dict>
		<key>com.apple.mobile.lockdown.remote.trusted</key>
		<dict>
			<key>ExposedToUntrustedDevices</key>
			<false/>
			<key>RequireEntitlement</key>
			<string>com.apple.mobile.lockdown.remote.trusted</string>
			<key>ServiceProperties</key>
			<dict>
				<key>ServiceVersion</key>
				<integer>1</integer>
				<key>UsesRemoteXPC</key>
				<false/>
			</dict>
		</dict>
		<key>com.apple.mobile.lockdown.remote.untrusted</key>
		<dict>
			<key>ExposedToUntrustedDevices</key>
			<true/>
			<key>RequireEntitlement</key>
			<string>com.apple.mobile.lockdown.remote.untrusted</string>
			<key>ServiceProperties</key>
			<dict>
				<key>ServiceVersion</key>
				<integer>1</integer>
				<key>UsesRemoteXPC</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>ServiceIPC</key>
	<false/>
	<key>Sockets</key>
	<dict>
		<key>MyIPv4Socket</key>
		<dict>
			<key>SockFamily</key>
			<string>IPv4</string>
			<key>SockServiceName</key>
			<string>62078</string>
		</dict>
		<key>MyIPv6Socket</key>
		<dict>
			<key>SockFamily</key>
			<string>IPv6</string>
			<key>SockServiceName</key>
			<string>62078</string>
		</dict>
		<key>MyUnixSocket</key>
		<dict>
			<key>SockFamily</key>
			<string>Unix</string>
			<key>SockPathMode</key>
			<integer>511</integer>
			<key>SockPathName</key>
			<string>/var/run/lockdown.sock</string>
		</dict>
	</dict>
	<key>_PanicOnCrash</key>
	<dict>
		<key>AutomationOnly</key>
		<true/>
		<key>PanicOnConsecutiveCrash</key>
		<true/>
	</dict>
</dict>
</plist>
