# beta/iOS/System/Library/LaunchDaemons/com.apple.dataaccess.dataaccessd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.dataaccess.dataaccessd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>NoteContextDarwinNotificationWithLoggedChanges</key>
			<dict>
				<key>Notification</key>
				<string>NoteContextDarwinNotificationWithLoggedChanges</string>
			</dict>
			<key>PCPreferencesDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>PCPreferencesDidChangeNotification</string>
			</dict>
			<key>_CalDatabaseChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>_CalDatabaseChangedNotification</string>
			</dict>
			<key>__ABDataBaseChangedByOtherProcessNotification</key>
			<dict>
				<key>Notification</key>
				<string>__ABDataBaseChangedByOtherProcessNotification</string>
			</dict>
			<key>com.apple.calendar.database.preference.notification.kCalPreferredDaysToSyncKey</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.calendar.database.preference.notification.kCalPreferredDaysToSyncKey</string>
			</dict>
			<key>com.apple.dataaccess.checkHolidayCalendarAccount</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.dataaccess.checkHolidayCalendarAccount</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.dataaccess.connected</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RequireBatteryLevel</key>
				<integer>15</integer>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequireScreenSleep</key>
				<false/>
				<key>RequiresBuddyComplete</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.dataaccess.dataaccessd</key>
		<true/>
		<key>com.apple.dataaccess.dataaccessd.aps</key>
		<true/>
		<key>com.apple.dataaccess.dataaccessd.aps.demo</key>
		<true/>
		<key>com.apple.dataaccess.dataaccessd.aps.dev</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/DataAccess.framework/Support/dataaccessd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
