# beta/iOS/System/Library/LaunchDaemons/com.apple.mobile.insecure_notification_proxy.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.mobile.insecure_notification_proxy</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mobile.insecure_notification_proxy</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/notification_proxy</string>
		<string>-i</string>
	</array>
	<key>RemoteServices</key>
	<dict>
		<key>com.apple.mobile.insecure_notification_proxy.remote</key>
		<dict>
			<key>ExposedToUntrustedDevices</key>
			<true/>
			<key>RequireEntitlement</key>
			<string>com.apple.mobile.insecure_notification_proxy.remote</string>
			<key>ServiceProperties</key>
			<dict>
				<key>ServiceVersion</key>
				<integer>1</integer>
				<key>UsesRemoteXPC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
