# beta/iOS/System/Library/LaunchDaemons/com.apple.asd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.asd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>PKPassLibraryDidAddPassNotification</key>
			<dict>
				<key>Name</key>
				<string>PKPassLibraryDidAddPassNotification</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.asd.odp</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>GracePeriod</key>
				<integer>43200</integer>
				<key>Interval</key>
				<integer>57600</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.asd.odp_a</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>43200</integer>
				<key>Interval</key>
				<integer>57600</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequireScreenSleep</key>
				<false/>
				<key>RequiresClassA</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.asdPushPort</key>
		<true/>
		<key>com.apple.asd.scoring</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/asd</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
