# beta/iOS/System/Library/LaunchDaemons/com.apple.mobile.storage_mounter_proxy.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.mobile.storage_mounter_proxy</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mobile.storage_mounter_proxy</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/mobile_storage_proxy</string>
	</array>
	<key>RemoteServices</key>
	<dict>
		<key>com.apple.mobile.storage_mounter_proxy.bridge</key>
		<dict>
			<key>ExposedToUntrustedDevices</key>
			<array>
				<string>HasInternalSecurityPolicies</string>
				<string>PreSetup</string>
			</array>
			<key>LimitExposedFromDeviceType</key>
			<array>
				<string>compute-node</string>
			</array>
			<key>RequireEntitlement</key>
			<string>com.apple.private.mobile_storage.remote.allowedSPI</string>
			<key>ServiceProperties</key>
			<dict>
				<key>ServiceVersion</key>
				<integer>1</integer>
				<key>UsesRemoteXPC</key>
				<true/>
			</dict>
		</dict>
	</dict>
</dict>
</plist>
