# beta/iOS/System/Library/LaunchDaemons/com.apple.ClarityBoard.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EmbeddedHomeScreen</key>
	<true/>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<false/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.ClarityBoard</string>
	<key>MachServices</key>
	<dict>
		<key>PurpleSystemAppPort</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.DragUI.druid.system</key>
		<true/>
		<key>com.apple.SBUserNotification</key>
		<true/>
		<key>com.apple.UIKit.KeyboardManagement.hosted</key>
		<true/>
		<key>com.apple.UIKit.statusbarserver</key>
		<true/>
		<key>com.apple.accessibility.AXClarityBoardServer</key>
		<true/>
		<key>com.apple.aps.usernotifications</key>
		<true/>
		<key>com.apple.aps.usernotifications.dev</key>
		<true/>
		<key>com.apple.bulletinboard.dataproviderconnection</key>
		<true/>
		<key>com.apple.bulletinboard.observerconnection</key>
		<true/>
		<key>com.apple.bulletinboard.serverconduitconnection</key>
		<true/>
		<key>com.apple.bulletinboard.settingsconnection</key>
		<true/>
		<key>com.apple.frontboard.systemappservices</key>
		<true/>
		<key>com.apple.frontboard.workspace</key>
		<true/>
		<key>com.apple.siri.activation</key>
		<true/>
		<key>com.apple.siri.activation.service</key>
		<true/>
		<key>com.apple.systemstatus</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.systemstatus.activityattribution</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.systemstatus.publisher</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.usernotifications.remotenotificationservice</key>
		<true/>
		<key>com.apple.usernotifications.systemservice</key>
		<true/>
		<key>com.apple.usernotifications.usernotificationservice</key>
		<true/>
		<key>com.apple.usernotifications.usernotificationsettingsservice</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/ClarityBoard.app/ClarityBoard</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UIApplicationInitializationContextFactoryClass</key>
	<string>ClarityBoard</string>
	<key>UserName</key>
	<string>mobile</string>
	<key>_LimitLoadToClarityMode</key>
	<true/>
</dict>
</plist>
