# beta/iOS/System/Library/LaunchDaemons/com.apple.cameracaptured.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<dict>
		<key>#Else</key>
		<false/>
		<key>#IfFeatureFlagEnabled</key>
		<string>CameraCapture/cameracaptured</string>
		<key>#Then</key>
		<true/>
	</dict>
	<key>Label</key>
	<string>com.apple.cameracaptured</string>
	<key>LimitLoadFromHardware</key>
	<dict>
		<key>machine</key>
		<array>
			<string>AudioAccessory1,1</string>
			<string>AudioAccessory1,2</string>
			<string>AudioAccessory5,1</string>
			<string>AudioAccessory6,1</string>
		</array>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>#IfFeatureFlagEnabled</key>
		<string>CameraCapture/cameracaptured</string>
		<key>#Then</key>
		<dict>
			<key>com.apple.coremedia.cameraviewfinder</key>
			<true/>
			<key>com.apple.coremedia.capture.deferred.containermanager.xpc</key>
			<true/>
			<key>com.apple.coremedia.capturesession</key>
			<true/>
			<key>com.apple.coremedia.capturesource</key>
			<true/>
			<key>com.apple.coremedia.externalstoragedevicemanager.xpc</key>
			<true/>
			<key>com.apple.coremedia.figvirtualcapturecard.xpc</key>
			<true/>
			<key>com.apple.coremedia.flashlight</key>
			<true/>
			<key>com.apple.coremedia.sandboxserver.xpc</key>
			<true/>
			<key>com.apple.mediasafetynet.exceptions.cam</key>
			<true/>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/cameracaptured</string>
	</array>
	<key>PublishesEvents</key>
	<dict>
		<key>#IfFeatureFlagEnabled</key>
		<string>CameraCapture/cameracaptured</string>
		<key>#Then</key>
		<array>
			<string>com.apple.cmio.registerassistantservice.system-extensions.matching</string>
		</array>
	</dict>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>UserName</key>
	<string>mobile</string>
	<key>_Conclave</key>
	<string>com.apple.cameracaptured.conclave</string>
</dict>
</plist>
