# beta/iOS/System/Library/LaunchDaemons/com.apple.routined.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.routined</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.accounts.changes</key>
		<dict>
			<key>AppleAccountAddedOrDeleted</key>
			<dict>
				<key>com.apple.account.AppleAccount</key>
				<dict>
					<key>ChangeTypes</key>
					<array>
						<integer>1</integer>
						<integer>3</integer>
					</array>
				</dict>
			</dict>
		</dict>
		<key>com.apple.alarm</key>
		<true/>
		<key>com.apple.commcenter.DataSettingsChangedNotification</key>
		<true/>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CNContactStoreDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>CNContactStoreDidChangeNotification</string>
			</dict>
			<key>CNContactStoreMeContactDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>CNContactStoreMeContactDidChangeNotification</string>
			</dict>
			<key>HKHealthDaemonActiveWorkoutServersDidUpdateNotification</key>
			<dict>
				<key>Notification</key>
				<string>HKHealthDaemonActiveWorkoutServersDidUpdateNotification</string>
			</dict>
			<key>com.apple.SOSEngine.SOSNotifyContactsReasonCinnamon</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SOSEngine.SOSNotifyContactsReasonCinnamon</string>
			</dict>
			<key>com.apple.SOSEngine.SOSNotifyContactsReasonKappa</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SOSEngine.SOSNotifyContactsReasonKappa</string>
			</dict>
			<key>com.apple.SOSEngine.SOSNotifyContactsReasonMandrake</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SOSEngine.SOSNotifyContactsReasonMandrake</string>
			</dict>
			<key>com.apple.SOSEngine.SOSNotifyContactsReasonNewton</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SOSEngine.SOSNotifyContactsReasonNewton</string>
			</dict>
			<key>com.apple.SOSEngine.SOSNotifyContactsReasonSOSTrigger</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SOSEngine.SOSNotifyContactsReasonSOSTrigger</string>
			</dict>
			<key>com.apple.isp.backcamerapower</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.isp.backcamerapower</string>
			</dict>
			<key>com.apple.isp.frontcamerapower</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.isp.frontcamerapower</string>
			</dict>
			<key>com.apple.locationd.appreset</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.appreset</string>
			</dict>
			<key>com.apple.locationd.gathering.bufferedDevicesReceived</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.gathering.bufferedDevicesReceived</string>
			</dict>
			<key>com.apple.locationd.gathering.came_up</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.gathering.came_up</string>
			</dict>
			<key>com.apple.locationd.vehicle.connected</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.vehicle.connected</string>
			</dict>
			<key>com.apple.locationd.vehicle.disconnected</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.vehicle.disconnected</string>
			</dict>
			<key>com.apple.locationd.vehicle.exit</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.vehicle.exit</string>
			</dict>
			<key>com.apple.mobile.keybagd.lock_status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.lock_status</string>
			</dict>
			<key>com.apple.siri.cloud.synch.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.cloud.synch.changed</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.routined</key>
		<true/>
		<key>com.apple.private.alloy.safetymonitor-idswake</key>
		<true/>
		<key>com.apple.private.alloy.safetymonitor.multiplex-idswake</key>
		<true/>
		<key>com.apple.private.alloy.safetymonitor.ownaccount-idswake</key>
		<true/>
		<key>com.apple.routined.event</key>
		<true/>
		<key>com.apple.routined.internal</key>
		<true/>
		<key>com.apple.routined.registration</key>
		<true/>
		<key>com.apple.routined.safetyMonitor</key>
		<true/>
		<key>com.apple.routined.store.cache</key>
		<true/>
		<key>com.apple.routined.store.cloud</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/routined</string>
		<string>LAUNCHED_BY_LAUNCHD</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
