# beta/iOS/System/Library/LaunchDaemons/com.apple.mobile_diagnostics_relay.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.mobile_diagnostics_relay</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.usbhub.attach.4110</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IOUSBHostDevice</string>
				<key>idProduct</key>
				<integer>4110</integer>
				<key>idVendor</key>
				<integer>1452</integer>
			</dict>
			<key>com.apple.usbhub.attach.4111</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IOUSBHostDevice</string>
				<key>idProduct</key>
				<integer>4111</integer>
				<key>idVendor</key>
				<integer>1452</integer>
			</dict>
			<key>com.apple.usbhub.attach.4119</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IOUSBHostDevice</string>
				<key>idProduct</key>
				<integer>4119</integer>
				<key>idVendor</key>
				<integer>1452</integer>
			</dict>
			<key>com.apple.usbhub.attach.4120</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IOUSBHostDevice</string>
				<key>idProduct</key>
				<integer>4120</integer>
				<key>idVendor</key>
				<integer>1452</integer>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mobile.diagnostics_relay.lockdown</key>
		<true/>
		<key>com.apple.mobile_diagnostics_relay</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/mobile_diagnostics_relay</string>
	</array>
</dict>
</plist>
