# beta/iOS/System/Library/LaunchDaemons/com.apple.safetyalertsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.safetyalertsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.CTTelephonyCenter</key>
		<dict>
			<key>kCTEventSMSMessageCellBroadcastStringGeofencingStarted</key>
			<dict>
				<key>Name</key>
				<string>kCTEventSMSMessageCellBroadcastStringGeofencingStarted</string>
			</dict>
			<key>kCTSMSCellBroadcastStringReceivedNotification</key>
			<dict>
				<key>Name</key>
				<string>kCTSMSCellBroadcastStringReceivedNotification</string>
			</dict>
			<key>kCTSMSCellBroadcastTelemetryNotification</key>
			<dict>
				<key>Name</key>
				<string>kCTSMSCellBroadcastTelemetryNotification</string>
			</dict>
		</dict>
		<key>com.apple.bluetooth.discovery</key>
		<dict>
			<key>SafetyAlertsBLERelay</key>
			<dict>
				<key>bleScanRate</key>
				<string>Low</string>
				<key>bleScanRateScreenOff</key>
				<string>Background</string>
				<key>discoveryFlags</key>
				<array>
					<string>ScreenOff</string>
				</array>
				<key>discoveryTypes</key>
				<array>
					<string>SafetyAlerts</string>
				</array>
				<key>label</key>
				<string>SafetyAlertsBLERelay</string>
				<key>xpcReportCompleteDevice</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>kCTSMSCellBroadcastConfigChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>kCTSMSCellBroadcastConfigChangedNotification</string>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToHardware</key>
	<dict>
		<key>targettype</key>
		<array>
			<string>N228s</string>
			<string>N244b</string>
			<string>N244s</string>
			<string>N230</string>
			<string>V57</string>
			<string>V54</string>
			<string>V53</string>
			<string>D23</string>
			<string>D16</string>
			<string>D17</string>
			<string>D20</string>
			<string>D201</string>
			<string>D21</string>
			<string>D211</string>
			<string>D22</string>
			<string>D221</string>
			<string>D27</string>
			<string>D28</string>
			<string>D321</string>
			<string>D331</string>
			<string>D331p</string>
			<string>D421</string>
			<string>D431</string>
			<string>D49</string>
			<string>D52g</string>
			<string>D53g</string>
			<string>D53p</string>
			<string>D54p</string>
			<string>D63</string>
			<string>D64</string>
			<string>D73</string>
			<string>D74</string>
			<string>D79</string>
			<string>D37</string>
			<string>D38</string>
			<string>D83</string>
			<string>D84</string>
			<string>D47</string>
			<string>D48</string>
			<string>D93</string>
			<string>D94</string>
			<string>N104</string>
			<string>N841</string>
			<string>V59</string>
			<string>N141b</string>
			<string>N141s</string>
			<string>N142b</string>
			<string>N142s</string>
			<string>N146b</string>
			<string>N146s</string>
			<string>N149b</string>
			<string>N149s</string>
			<string>N158b</string>
			<string>N158s</string>
			<string>N188b</string>
			<string>N188s</string>
			<string>N198b</string>
			<string>N198s</string>
			<string>N199</string>
			<string>N208b</string>
			<string>N208s</string>
			<string>N210</string>
			<string>N218s</string>
			<string>N218b</string>
			<string>N143b</string>
			<string>N143s</string>
			<string>N157b</string>
			<string>N157s</string>
			<string>N187b</string>
			<string>N187s</string>
			<string>N197b</string>
			<string>N197s</string>
			<string>N207b</string>
			<string>N207s</string>
			<string>N217b</string>
			<string>N217s</string>
			<string>N227b</string>
			<string>N227s</string>
			<string>N228b</string>
			<string>N243b</string>
			<string>N243s</string>
		</array>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.safetyalerts</key>
		<true/>
		<key>com.apple.private.alloy.safetyalerts-idswake</key>
		<true/>
		<key>com.apple.safetyalerts</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.safetyalerts</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/safetyalertsd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
