# beta/iOS/System/Library/LaunchDaemons/com.apple.inputanalyticsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.inputanalyticsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.inputanalyticsd.daily</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>82800</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.inputanalyticsd</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/inputanalyticsd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
