# beta/iOS/System/Library/LaunchDaemons/com.apple.aeroml.intentrecommend.intentrecommendd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.aeroml.intentrecommend.intentrecommendd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.biome.compute.publisher</key>
		<dict>
			<key>com.apple.intentrecommendd.App.Intent</key>
			<dict>
				<key>client</key>
				<string>com.apple.aeroml.intentrecommend.intentrecommendd</string>
				<key>graph</key>
				<dict>
					<key>keyPath</key>
					<string>eventBody.intentClass</string>
					<key>name</key>
					<string>filter</string>
					<key>upstream</key>
					<dict>
						<key>name</key>
						<string>streamPublisher</string>
						<key>streamIdentifier</key>
						<string>App.Intent</string>
						<key>version</key>
						<integer>1</integer>
					</dict>
					<key>value</key>
					<string>INPlayMediaIntent</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.aeroml.intentrecommendd.App.Intent</string>
			</dict>
			<key>com.apple.intentrecommendd.App.Intents.Transcript</key>
			<dict>
				<key>client</key>
				<string>com.apple.aeroml.intentrecommend.intentrecommendd</string>
				<key>graph</key>
				<dict>
					<key>dsl&lt;@BMDSLFilter&gt;</key>
					<dict>
						<key>comparison</key>
						<integer>3</integer>
						<key>keyPath</key>
						<string>eventBody.bundleID</string>
						<key>upstream</key>
						<dict>
							<key>name</key>
							<string>streamPublisher</string>
							<key>streamIdentifier</key>
							<string>App.Intents.Transcript</string>
							<key>version</key>
							<integer>1</integer>
						</dict>
						<key>value</key>
						<array>
							<string>com.apple.podcasts</string>
							<string>com.apple.NanoMusic</string>
							<string>com.apple.Music</string>
						</array>
					</dict>
				</dict>
				<key>identifier</key>
				<string>com.apple.intentrecommendd.App.Intents.Transcript</string>
			</dict>
		</dict>
		<key>com.apple.private.intelligenceplatform.use-cases</key>
		<dict>
			<key>IntentRecommendEvent</key>
			<dict>
				<key>Streams</key>
				<array>
					<string>App.Intent</string>
					<string>App.Intents.Transcript</string>
				</array>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict/>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aeroml.intentrecommend.mediasuggester</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/IntentRecommendRuntime.framework/intentrecommendd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
