# beta/iOS/System/Library/LaunchDaemons/com.apple.fusiond.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.fusiond</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.fusiond</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/fusiond</string>
	<key>RemoteServices</key>
	<dict>
		<key>com.apple.fusion.remote.service</key>
		<dict>
			<key>ExposedToUntrustedDevices</key>
			<true/>
			<key>RequireEntitlement</key>
			<string>com.apple.fusion.remote.service</string>
			<key>ServiceProperties</key>
			<dict>
				<key>ServiceVersion</key>
				<integer>1</integer>
				<key>UsesRemoteXPC</key>
				<true/>
			</dict>
		</dict>
	</dict>
</dict>
</plist>
