# beta/iOS/System/Library/LaunchDaemons/com.apple.financed.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.financed</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.financed.bankconnect.lab</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>259200</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.financed.instrumentation</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.financed.maintenance</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.financed.orders.reboot</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.financed.post-install</key>
			<dict>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.financed.refresh.institutions</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>2592000</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RequiresUserInactivity</key>
				<false/>
			</dict>
		</dict>
		<key>com.apple.biome.compute.publisher</key>
		<dict>
			<key>com.apple.financed.biome.extractedOrder</key>
			<dict>
				<key>client</key>
				<string>com.apple.financed</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>WalletPaymentsCommerce.FoundIn.ExtractedOrder</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.financed.biome.extractedOrder</string>
			</dict>
		</dict>
		<key>com.apple.cascade.xpc_event.setChange</key>
		<dict>
			<key>com.apple.financed.cascade.extractedOrder</key>
			<dict/>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.financed.bankconnect.reboot</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequiresClassC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.financed</key>
		<true/>
		<key>com.apple.aps.financed.store</key>
		<true/>
		<key>com.apple.financed.service.bankconnect</key>
		<true/>
		<key>com.apple.financed.service.coredatastore</key>
		<true/>
		<key>com.apple.financed.service.financestore</key>
		<true/>
		<key>com.apple.financed.service.store</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/FinanceKit.framework/financed</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
