# beta/iOS/System/Library/LaunchDaemons/com.apple.SecureBackupDaemon.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.SecureBackupDaemon</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.sbd.kvstorechange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sbd.kvstorechange</string>
			</dict>
			<key>com.apple.security.itembackup</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.security.itembackup</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.SecureBackupDaemon</key>
		<true/>
		<key>com.apple.SecureBackupDaemon.concurrent</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CloudServices.framework/Helpers/com.apple.sbd</string>
	</array>
</dict>
</plist>
