# beta/iOS/System/Library/LaunchDaemons/com.apple.mobilerepaird.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.mobilerepaird</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.accessories.hwComponent.isTrustedForUI.failure.battery</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>InternalComponent</key>
					<true/>
					<key>isTrustedForUI</key>
					<false/>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleAuthCPRelay</string>
			</dict>
			<key>com.apple.accessories.hwComponent.isTrustedForUI.failure.display</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>InternalComponent</key>
					<true/>
					<key>isTrustedForUI</key>
					<false/>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleAuthCPAID</string>
			</dict>
			<key>com.apple.accessories.hwComponent.isTrustedForUI.pass.battery</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>InternalComponent</key>
					<true/>
					<key>isTrustedForUI</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleAuthCPRelay</string>
			</dict>
			<key>com.apple.accessories.hwComponent.isTrustedForUI.pass.display</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>InternalComponent</key>
					<true/>
					<key>isTrustedForUI</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleAuthCPAID</string>
			</dict>
			<key>com.apple.audio.hw.failure.codec</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>AudioHardwareError</key>
					<string>Codec</string>
				</dict>
			</dict>
			<key>com.apple.audio.hw.failure.speaker</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>AudioHardwareError</key>
					<string>Speaker</string>
				</dict>
			</dict>
			<key>com.apple.camera.fdr.failure.h13</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>CmClValidationStatus</key>
					<string>Fail</string>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleH13CamIn</string>
			</dict>
			<key>com.apple.camera.fdr.failure.h16</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>CmClValidationStatus</key>
					<string>Fail</string>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleH16CamIn</string>
			</dict>
			<key>com.apple.camera.fdr.invalid.h13</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>CmClValidationStatus</key>
					<string>Invalid</string>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleH13CamIn</string>
			</dict>
			<key>com.apple.camera.fdr.invalid.h16</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>CmClValidationStatus</key>
					<string>Invalid</string>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleH16CamIn</string>
			</dict>
			<key>com.apple.camera.fdr.pass.h13</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>CmClValidationStatus</key>
					<string>Pass</string>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleH13CamIn</string>
			</dict>
			<key>com.apple.camera.fdr.pass.h16</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>CmClValidationStatus</key>
					<string>Pass</string>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleH16CamIn</string>
			</dict>
			<key>com.apple.camera.hw.failure.h10</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>BackCameraExpected</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleH10CamIn</string>
			</dict>
			<key>com.apple.camera.hw.failure.h13</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>BackCameraExpected</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleH13CamIn</string>
			</dict>
			<key>com.apple.camera.hw.failure.h16</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>BackCameraExpected</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleH16CamIn</string>
			</dict>
			<key>com.apple.camera.hw.failure.h9</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>BackCameraExpected</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleH9CamIn</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AppleLanguagePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleLanguagePreferencesChangedNotification</string>
			</dict>
			<key>com.apple.mobilerepair.refreshui</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobilerepair.refreshui</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.mobilerepaird.cache_clean</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
			</dict>
			<key>com.apple.mobilerepaird.daily_stats</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.mobilerepaird.post_install</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mobilerepaird</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/mobilerepaird</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
