# beta/iOS/System/Library/LaunchDaemons/com.apple.tipsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.tipsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.coreduetcontext.client_event_stream</key>
		<dict>
			<key>notification</key>
			<dict>
				<key>client-identifier</key>
				<string>com.apple.tipsd</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.homed.speakersConfiguredChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.speakersConfiguredChanged</string>
			</dict>
			<key>com.apple.language.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.language.changed</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.purplebuddy.setupexited</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.purplebuddy.setupexited</string>
			</dict>
		</dict>
		<key>com.apple.nsurlsessiond-events</key>
		<dict>
			<key>LaunchOnDemand</key>
			<dict>
				<key>ClientIdentifier</key>
				<string>com.apple.tipsd</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.tipsd.analytics-maintenance</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.tipsd.postInstall-utility</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>NetworkTransferDirection</key>
				<string>Download</string>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.tipsd.update-content</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>7200</integer>
				<key>Interval</key>
				<integer>21600</integer>
				<key>NetworkTransferDirection</key>
				<string>Download</string>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
				<key>RequireScreenSleep</key>
				<false/>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.tipsd</key>
		<true/>
		<key>com.apple.tipsd.assistant</key>
		<true/>
		<key>com.apple.tipsd.supportFlow.analytics</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/tipsd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
