# beta/iOS/System/Library/LaunchDaemons/com.apple.companiond.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.companiond</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bluetooth.connections</key>
		<dict>
			<key>CompanionServices</key>
			<dict>
				<key>label</key>
				<string>CompanionServices</string>
				<key>psm</key>
				<integer>136</integer>
			</dict>
		</dict>
		<key>com.apple.bluetooth.discovery</key>
		<dict>
			<key>PhotoSetup</key>
			<dict>
				<key>authFlags</key>
				<array>
					<string>SameAccount</string>
				</array>
				<key>bleRSSIThresholdHint</key>
				<integer>-70</integer>
				<key>discoveryTypes</key>
				<array>
					<string>PhotoSetup</string>
				</array>
				<key>label</key>
				<string>PhotoSetup</string>
				<key>xpcReportCompleteDevice</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>CPSRemoteLLM</key>
			<dict>
				<key>serviceType</key>
				<string>CPSRemoteLLM</string>
				<key>type</key>
				<string>server</string>
			</dict>
			<key>CompanionAuthenticationDiscovery</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.CompanionAuthentication</string>
				<key>type</key>
				<string>discovery</string>
			</dict>
			<key>LegacyAppSignInDiscovery</key>
			<dict>
				<key>serviceType</key>
				<string>RPUserNotificationAppSignIn</string>
				<key>type</key>
				<string>discovery</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.CompanionServices</key>
		<true/>
		<key>com.apple.DistributedTimers</key>
		<true/>
		<key>com.apple.companiond.xpc</key>
		<true/>
		<key>com.apple.private.alloy.companion.auth-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.CompanionNotifications</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/companiond</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
