# beta/iOS/System/Library/LaunchDaemons/com.apple.siriactionsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MallocStackLoggingDirectory</key>
		<string>/private/var/tmp/com.apple.siriactionsd</string>
	</dict>
	<key>Label</key>
	<string>com.apple.siriactionsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.appprotection.preferences.change.event</key>
		<dict>
			<key>Hidden App Bundle IDs Subscription</key>
			<dict>
				<key>hiddenEvent</key>
				<true/>
				<key>unhiddenEvent</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.biome.compute.publisher</key>
		<dict>
			<key>com.apple.siriactionsd.UserFocusComputedMode</key>
			<dict>
				<key>client</key>
				<string>com.apple.siriactionsd</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>UserFocusComputedMode</string>
					<key>version</key>
					<string>1</string>
				</dict>
				<key>identifier</key>
				<string>com.apple.siriactionsd.UserFocusComputedMode</string>
			</dict>
		</dict>
		<key>com.apple.cascade.xpc_event.setChange</key>
		<dict>
			<key>com.apple.toolkit.sync.pull</key>
			<dict/>
		</dict>
		<key>com.apple.coreduetcontext.client_event_stream</key>
		<dict>
			<key>notification</key>
			<dict>
				<key>client-identifier</key>
				<string>com.apple.siriactionsd.contextstore-registration</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>ApplicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>ApplicationStateChanged</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationStateChanged</string>
			</dict>
			<key>ApplicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
			<key>CoreDataDatabaseChanged</key>
			<dict>
				<key>Name</key>
				<string>com.apple.shortcuts.WFCoreDataDatabaseContextDidSaveNotification</string>
			</dict>
			<key>LNAppShortcutsChangedNotification</key>
			<dict>
				<key>Name</key>
				<string>LNAppShortcutsChangedNotification</string>
			</dict>
			<key>LNMetadataChangedNotification</key>
			<dict>
				<key>Name</key>
				<string>LNMetadataChangedNotification</string>
			</dict>
			<key>TKToolkitDatabaseChangedNotification</key>
			<dict>
				<key>Name</key>
				<string>TKToolkitDatabaseChangedNotification</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CNFavoritesChangedExternallyNotification</key>
			<dict>
				<key>Notification</key>
				<string>CNFavoritesChangedExternallyNotification</string>
			</dict>
			<key>UIAccessibilityInvertColorsChanged</key>
			<dict>
				<key>Notification</key>
				<string>UIAccessibilityInvertColorsChanged</string>
			</dict>
			<key>com.apple.accessibility.classic.wob.status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.classic.wob.status</string>
			</dict>
			<key>com.apple.accessibility.commandandcontrol.status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.commandandcontrol.status</string>
			</dict>
			<key>com.apple.accessibility.enhance.background.contrast.status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.enhance.background.contrast.status</string>
			</dict>
			<key>com.apple.accessibility.pointer.increased.contrast</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.pointer.increased.contrast</string>
			</dict>
			<key>com.apple.accessibility.reduce.motion.status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.reduce.motion.status</string>
			</dict>
			<key>com.apple.accessibility.reduce.white.point</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.reduce.white.point</string>
			</dict>
			<key>com.apple.accessibility.voiceovertouch.status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.voiceovertouch.status</string>
			</dict>
			<key>com.apple.accessibility.zoomtouch.status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.accessibility.zoomtouch.status</string>
			</dict>
			<key>com.apple.appprotection.change.hidden</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appprotection.change.hidden</string>
			</dict>
			<key>com.apple.callhistory.notification.calls-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.callhistory.notification.calls-changed</string>
			</dict>
			<key>com.apple.commcenter.DataSettingsChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.commcenter.DataSettingsChangedNotification</string>
			</dict>
			<key>com.apple.language.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.language.changed</string>
			</dict>
			<key>com.apple.mediaaccessibility.displayFilterSettingsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mediaaccessibility.displayFilterSettingsChanged</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.nanoregistry.paireddevicedidchangeversion</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.paireddevicedidchangeversion</string>
			</dict>
			<key>com.apple.shortcuts.daemon-wakeup-request</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.shortcuts.daemon-wakeup-request</string>
			</dict>
			<key>com.apple.siri.ShortcutsCloudKitAccountAddedNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.ShortcutsCloudKitAccountAddedNotification</string>
			</dict>
			<key>com.apple.siri.ShortcutsCloudKitAccountModifiedNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.ShortcutsCloudKitAccountModifiedNotification</string>
			</dict>
			<key>com.apple.spotlightui.prefschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.spotlightui.prefschanged</string>
			</dict>
			<key>com.apple.system.clock_set</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.clock_set</string>
			</dict>
			<key>com.apple.system.timezone</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.timezone</string>
			</dict>
			<key>com.apple.workflow.remotewidgets</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.workflow.remotewidgets</string>
			</dict>
			<key>kAFPreferencesDidChangeDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>kAFPreferencesDidChangeDarwinNotification</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>remote-widgets-connection</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.workflow.remotewidgets</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.siriactionsd.RemoteQuarantineUpdate</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>DuetRelatedApplications</key>
				<array>
					<string>com.apple.shortcuts</string>
				</array>
				<key>Interval</key>
				<integer>86400</integer>
				<key>NetworkTransferDirection</key>
				<string>Download</string>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>3600</integer>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.siriactionsd.ShortcutsDailyMetricCheckIn</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>DuetRelatedApplications</key>
				<array>
					<string>com.apple.shortcuts</string>
				</array>
				<key>Interval</key>
				<integer>28800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<real>28800</real>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.siriactionsd.ShortcutsWeeklyMetricCheckIn</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>DuetRelatedApplications</key>
				<array>
					<string>com.apple.shortcuts</string>
				</array>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<integer>86400</integer>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.siriactionsd.TriggersEnabledBootNotification</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>0</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequiresClassC</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.siriactionsd</key>
		<true/>
		<key>com.apple.private.alloy.shortcuts-idswake</key>
		<true/>
		<key>com.apple.private.alloy.siri.voiceshortcuts-idswake</key>
		<true/>
		<key>com.apple.shortcuts.dialogpresentation</key>
		<true/>
		<key>com.apple.siri.VoiceShortcuts.xpc</key>
		<true/>
		<key>com.apple.siriactionsd.pairedsync</key>
		<true/>
		<key>com.apple.siriactionsd.xpc</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.shortcuts</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.shortcuts.AutomationLockscreenNotifications</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/VoiceShortcuts.framework/Support/siriactionsd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
