# beta/iOS/System/Library/LaunchDaemons/com.apple.identityservicesd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>NSRunningFromLaunchd</key>
		<string>1</string>
		<key>OS_ACTIVITY_ENABLE_DYNAMIC</key>
		<string>1</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.apple.identityservicesd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.identityservicesd.aps</key>
		<true/>
		<key>com.apple.identityservicesd.embedded.auth</key>
		<true/>
		<key>com.apple.identityservicesd.idquery.embedded.auth</key>
		<true/>
		<key>com.apple.identityservicesd.nsxpc</key>
		<true/>
		<key>com.apple.identityservicesd.nsxpc.auth</key>
		<true/>
		<key>com.apple.identityservicesd.pds</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.identityservicesd.firewall</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/IDS.framework/identityservicesd.app/identityservicesd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
	<key>WorkingDirectory</key>
	<string>/tmp</string>
</dict>
</plist>
