# beta/iOS/System/Library/LaunchDaemons/com.apple.SyncAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.SyncAgent</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.lockdown.SyncAgent</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/SyncAgent</string>
		<string>--lockdown</string>
		<string>--oneshot</string>
		<string>-v</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
