# beta/iOS/System/Library/LaunchDaemons/com.apple.sysdiagnose.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EfficiencyMode</key>
	<string>Efficient</string>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.sysdiagnose</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.private.alloy.sysdiagnose-idswake</key>
		<true/>
		<key>com.apple.sysdiagnose.CacheDelete</key>
		<true/>
		<key>com.apple.sysdiagnose.kernel.ipc</key>
		<dict>
			<key>HostSpecialPort</key>
			<integer>23</integer>
		</dict>
		<key>com.apple.sysdiagnose.service.xpc</key>
		<dict/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/sysdiagnosed</string>
	</array>
	<key>RemoteServices</key>
	<dict>
		<key>com.apple.sysdiagnose.remote</key>
		<dict>
			<key>ExposedToUntrustedDevices</key>
			<array>
				<string>HasInternalSecurityPolicies</string>
			</array>
			<key>RequireEntitlement</key>
			<string>com.apple.private.sysdiagnose.remote</string>
			<key>ServiceProperties</key>
			<dict>
				<key>UsesRemoteXPC</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.sysdiagnose.remote.trusted</key>
		<dict>
			<key>RequireEntitlement</key>
			<string>com.apple.private.sysdiagnose.remote.trusted</string>
			<key>ServiceProperties</key>
			<dict>
				<key>UsesRemoteXPC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>_LimitLoadFromVariant</key>
	<string>IsDarwinOS</string>
</dict>
</plist>
