# diff/tvOS/System/Library/DeveloperModeLaunchDaemons/tmp407b17a4.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Group</key>
	<string>mobile</string>
	<key>Label</key>
	<string>com.apple.gputoolsserviced</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.gputools.events</key>
		<dict>
			<key>com.apple.gputools.events.TransportDaemonLaunched</key>
			<dict/>
		</dict>
	</dict>
	<key>LimitLoadToDeveloperMode</key>
	<true/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.gputools.service</key>
		<true/>
		<key>com.apple.gputools.service.simulator</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/usr/libexec/gputoolsserviced</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
