# beta/iOS/System/Library/NanoLaunchDaemons/com.apple.magicswitchd.companion.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/var/mobile/Library/DeviceRegistry/MagicSwitchKeepAlive</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.magicswitchd.companion</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>NRPairedDeviceRegistryPairedDeviceDidChangeCapabilitiesDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.paireddevicedidchangecapabilities</string>
			</dict>
			<key>NRPairedDeviceRegistryPairedDeviceDidChangeVersionDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.paireddevicedidchangeversion</string>
			</dict>
			<key>PSYDeviceSyncDidCompleteDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.pairedsync.syncDidComplete</string>
			</dict>
			<key>com.apple.NanoRegistry.MagicSwitchSettingsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.NanoRegistry.MagicSwitchSettingsChanged</string>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/magicswitchd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>TransactionTimeLimitEnabled</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
