# beta/System/Library/LaunchDaemons/com.apple.Maps.pushdaemon.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.Maps.mapspushd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AppleLanguagePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleLanguagePreferencesChangedNotification</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.Maps.mapspushd</key>
		<true/>
		<key>com.apple.Maps.mapspushd.geoservices</key>
		<true/>
		<key>com.apple.Maps.mapspushd.icloud</key>
		<true/>
		<key>com.apple.Maps.xpc.MediaAPI</key>
		<true/>
		<key>com.apple.Maps.xpc.SharedTrip</key>
		<true/>
		<key>com.apple.Maps.xpc.SharedTrip.Capabilities</key>
		<true/>
		<key>com.apple.Maps.xpc.connectionBroker.endpointProvider</key>
		<true/>
		<key>com.apple.Maps.xpc.connectionBroker.endpointRecorder</key>
		<true/>
		<key>com.apple.private.alloy.maps-idswake</key>
		<true/>
		<key>com.apple.private.alloy.maps.eta-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.Maps</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/MapsSupport.framework/mapspushd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/MapsSupport.framework/mapspushd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>SBMachServices</key>
	<array>
		<string>com.apple.Maps.xpc.connectionBroker.endpointRecorder</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
