# beta/System/Library/LaunchDaemons/com.apple.OTATaskingAgent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.OTATaskingAgent</string>
	<key>LimitLoadToSessionType</key>
	<dict>
		<key>#IfFeatureFlagEnabled</key>
		<string>UserManagement/SystemSessionD1</string>
		<key>#Then</key>
		<string>System</string>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.OTATaskingAgent</key>
		<dict/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/OTATaskingAgent</string>
		<string>server-init</string>
	</array>
</dict>
</plist>
