# beta/System/Library/LaunchDaemons/com.apple.eyereliefd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.eyereliefd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.ManagedSettings.effective-settings.changed</key>
		<dict>
			<key>com.apple.eyereliefd.managed-settings-changed</key>
			<dict>
				<key>SettingsGroups</key>
				<array>
					<string>eyeRelief</string>
				</array>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.eyereliefd.startup</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>120</integer>
				<key>Repeating</key>
				<false/>
				<key>RequiresClassA</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.eyereliefd</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/EyeRelief.framework/Resources/eyereliefd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
