# beta/System/Library/LaunchDaemons/com.apple.bird.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>DrainMessagesOnFailedInit</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.bird</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.bird.analytics-report</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>69120</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.bird.app-telemetry</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>172800</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>138240</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.bird.configuration-server-update</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>172800</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.bird.apps.install</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.bird.apps.uninstall</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.bird</key>
		<true/>
		<key>com.apple.bird.cache-delete</key>
		<true/>
		<key>com.apple.bird.push</key>
		<true/>
		<key>com.apple.bird.token</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.bird.usernotifications</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/iCloudDriveCore.framework/bird</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
