# beta/System/Library/LaunchDaemons/com.apple.ind.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamMemoryLimit</key>
		<real>5000</real>
	</dict>
	<key>Label</key>
	<string>com.apple.ind</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.dispatch.vfs</key>
		<dict>
			<key>Monitor Low Disk Conditions</key>
			<dict>
				<key>DesiredDisk</key>
				<true/>
				<key>LowDisk</key>
				<integer>60</integer>
				<key>NearLowDisk</key>
				<true/>
				<key>VeryLowDisk</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>DeviceDidPairEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
			<key>DeviceNameChangeEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.lockdown.device_name_changed</string>
			</dict>
			<key>FirstUnlockEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>LanguageChangedEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.language.changed</string>
			</dict>
			<key>PhotosOptimizeStateChangedEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobileslideshow.PLNotificationKeepOriginalsChanged</string>
			</dict>
			<key>PhotosiCPLStateChangedEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobileslideshow.ICPLStateChanged</string>
			</dict>
			<key>SimulateVFSAlmostFullEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.cloud.quota.simulate.vfs.almostfull</string>
			</dict>
			<key>SimulateVFSNotFullEvent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.cloud.quota.simulate.vfs.notfull</string>
			</dict>
			<key>SiriLanguageChangedEvent</key>
			<dict>
				<key>Notification</key>
				<string>kAFPreferencesDidChangeDarwinNotification</string>
			</dict>
			<key>refreshCFAvailability</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.os-eligibility-domain.change.greymatter</string>
			</dict>
			<key>refreshGMSAvailability</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.gms.availability.notification</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.ind.freshmint</key>
			<dict/>
			<key>com.apple.ind.heartbeat</key>
			<dict/>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.ind.and_callback</key>
		<true/>
		<key>com.apple.ind.aps</key>
		<true/>
		<key>com.apple.ind.xpc</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.iCloudSettingsNotification</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/iCloudNotification.framework/ind</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>3</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
