# beta/System/Library/LaunchDaemons/com.apple.memoryanalyticsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.memoryanalyticsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.da.tasking_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.da.tasking_changed</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.memoryanalyticsd.MTRandomMemgraphActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>ExpectedDuration</key>
				<integer>600</integer>
				<key>Interval</key>
				<integer>3600</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/memoryanalyticsd</string>
		<string>-xlatSf</string>
	</array>
</dict>
</plist>
