# beta/System/Library/LaunchDaemons/com.apple.purplebuddy.budd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.purplebuddy.budd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.managedconfiguration.effectivesettingschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.effectivesettingschanged</string>
			</dict>
			<key>com.apple.managedconfiguration.passcodechanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.passcodechanged</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.purplebuddy.backup</key>
			<dict>
				<key>DASData</key>
				<dict>
					<key>UserRequestedBackupActivity</key>
					<true/>
				</dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.purplebuddy.betaReEnrollment</key>
			<dict/>
			<key>com.apple.purplebuddy.finishsetup.firstreschedule</key>
			<dict/>
			<key>com.apple.purplebuddy.finishsetup.secondreschedule</key>
			<dict/>
			<key>com.apple.purplebuddy.finishsetup.wifiobserver</key>
			<dict/>
			<key>com.apple.purplebuddy.maintenance</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Interval</key>
				<integer>604800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.purplebuddy.budd.cloud.xpc</key>
		<true/>
		<key>com.apple.purplebuddy.budd.migration.source.xpc</key>
		<true/>
		<key>com.apple.purplebuddy.budd.proximity.source.xpc</key>
		<true/>
		<key>com.apple.purplebuddy.budd.proximity.target.xpc</key>
		<true/>
		<key>com.apple.purplebuddy.budd.proximityautomateddeviceenrollment.target.xpc</key>
		<true/>
		<key>com.apple.purplebuddy.budd.settings.xpc</key>
		<true/>
		<key>com.apple.purplebuddy.budd.xpc</key>
		<true/>
		<key>com.apple.purplebuddy.setup.budd.proximity.override.xpc</key>
		<true/>
		<key>com.apple.purplebuddy.setupctl.budd.proximity.override.service.xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/SetupAssistant.framework/budd</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>3</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
