# beta/System/Library/LaunchDaemons/com.apple.devicedataresetd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.devicedataresetd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.devicedatareset.DeviceDataResetObservationService.NonLaunching</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.devicedatareset.DeviceDataResetService</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/EmbeddedDataReset.framework/Support/devicedataresetd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
