# beta/System/Library/LaunchDaemons/com.apple.telephonyutilities.callservicesd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>NSRunningFromLaunchd</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.telephonyutilities.callservicesd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.CTTelephonyCenter</key>
		<dict>
			<key>com.apple.telephonyutilities.callservicesd.kCTSettingCallCapabilitiesChangedNotification</key>
			<dict>
				<key>Name</key>
				<string>kCTSettingCallCapabilitiesChangedNotification</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.callhistory.save.distributed.notification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.callhistory.save.distributed.notification</string>
			</dict>
			<key>com.apple.telephonyutilities.callservicesdaemon.app-installed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.telephonyutilities.callservicesdaemon.app-uninstalled</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.telephonyutilities.callservicesd.fakeincomingmessage</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.telephonyutilities.callservicesd.fakeincomingmessage</string>
			</dict>
			<key>com.apple.telephonyutilities.callservicesd.fakeoutgoingmessage</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.telephonyutilities.callservicesd.fakeoutgoingmessage</string>
			</dict>
			<key>com.apple.telephonyutilities.callservicesd.kFZACAppBundleIdentifierLaunchNotification</key>
			<dict>
				<key>Notification</key>
				<string>kFZACAppBundleIdentifierLaunchNotification</string>
			</dict>
			<key>com.apple.telephonyutilities.callservicesd.kFZVCAppBundleIdentifierLaunchNotification</key>
			<dict>
				<key>Notification</key>
				<string>kFZVCAppBundleIdentifierLaunchNotification</string>
			</dict>
			<key>com.apple.telephonyutilities.callservicesd.kFaceTimeChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>kFaceTimeChangedNotification</string>
			</dict>
			<key>com.apple.voicemail.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicemail.changed</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.callservicesd.ringtone.metrics.daily</key>
			<dict>
				<key>Delay</key>
				<real>86400</real>
				<key>GracePeriod</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>_SequenceNumber</key>
				<real>3</real>
			</dict>
			<key>com.apple.telephonyutilities.callservicesd.analytics.daily</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.telephonyutilities.callservicesd.analytics.weekly</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>604800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.telephonyutilities.callservicesd.calldirectorymaintenance.daily</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.telephonyutilities.callservicesd.csdvoipapplicationcontroller.daily</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<real>86400</real>
				<key>GracePeriod</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>_SequenceNumber</key>
				<real>3</real>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.callkit.callcontrollerhost</key>
		<true/>
		<key>com.apple.callkit.callsourcehost</key>
		<true/>
		<key>com.apple.callkit.networkextension.messagecontrollerhost</key>
		<true/>
		<key>com.apple.callkit.networkextension.voip</key>
		<true/>
		<key>com.apple.callkit.notificationserviceextension.voip</key>
		<true/>
		<key>com.apple.callkit.service</key>
		<true/>
		<key>com.apple.conversation.dialing</key>
		<true/>
		<key>com.apple.copresence.conversationmanagerhost</key>
		<true/>
		<key>com.apple.group-activities.conversationmanagerhost</key>
		<true/>
		<key>com.apple.private.alloy.dropin.communication-idswake</key>
		<true/>
		<key>com.apple.private.alloy.facetime.audio-idswake</key>
		<true/>
		<key>com.apple.private.alloy.facetime.multi-idswake</key>
		<true/>
		<key>com.apple.private.alloy.facetime.sync-idswake</key>
		<true/>
		<key>com.apple.private.alloy.facetime.video-idswake</key>
		<true/>
		<key>com.apple.private.alloy.phone.auth-idswake</key>
		<true/>
		<key>com.apple.private.alloy.phonecontinuity-idswake</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.callcapabilities</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.callhistorycontroller</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.callhistorymanager</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.callprovidermanager</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.callstatecontroller</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.conversationmanager</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.conversationprovidermanager</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.momentscontroller</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.neighborhood-activities</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.ptt</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.reportingcontroller</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.simulatedconversationcontroller</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.usernotificationprovider</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.voip</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.voip.push</key>
		<true/>
		<key>com.apple.telephonyutilities.callservicesdaemon.voip.push.development</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.facetime</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.mobilephone</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/TelephonyUtilities.framework/callservicesd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>_Conclave</key>
	<string>com.apple.callservicesd.conclave</string>
</dict>
</plist>
