# beta/System/Library/LaunchDaemons/com.apple.merchantd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.merchantd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.merchantd.serviceCleanup</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>100800</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>80640</integer>
				</dict>
				<key>RequiresProtectionClass</key>
				<string>A</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>TMRTCResetNotification</key>
			<dict>
				<key>Notification</key>
				<string>TMRTCResetNotification</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.merchantd.discovery</key>
		<true/>
		<key>com.apple.merchantd.identity</key>
		<true/>
		<key>com.apple.merchantd.internals</key>
		<true/>
		<key>com.apple.merchantd.storage</key>
		<true/>
		<key>com.apple.merchantd.transaction</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/merchantd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
