# beta/System/Library/LaunchDaemons/com.apple.carkitd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.carkitd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.biome.compute.publisher</key>
		<dict>
			<key>com.apple.internal.carkitd</key>
			<dict>
				<key>client</key>
				<string>com.apple.internal.carkitd</string>
				<key>graph</key>
				<dict>
					<key>keyPath</key>
					<string>eventBody.semanticType</string>
					<key>name</key>
					<string>filter</string>
					<key>upstream</key>
					<dict>
						<key>name</key>
						<string>streamPublisher</string>
						<key>streamIdentifier</key>
						<string>UserFocus.ComputedMode</string>
						<key>version</key>
						<integer>1</integer>
					</dict>
					<key>value</key>
					<integer>4</integer>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.focus.driving</string>
			</dict>
		</dict>
		<key>com.apple.bluetooth.connections</key>
		<dict>
			<key>CarPlay</key>
			<dict>
				<key>serviceUUIDs</key>
				<array>
					<string>93200A65-FDCF-472C-B1A8-8DB8A48DE8EA</string>
				</array>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.carkit.capabilities-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.carkit.capabilities-changed</string>
			</dict>
			<key>com.apple.carkit.carplay-attached</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.carkit.carplay-attached</string>
			</dict>
			<key>com.apple.locationd.vehicle.connected</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.vehicle.connected</string>
			</dict>
			<key>com.apple.locationd.vehicle.disconnected</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.vehicle.disconnected</string>
			</dict>
			<key>com.apple.vehiclePolicy.DNDMode</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.vehicular.changed.toVehicular</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.carkit.app.service</key>
		<true/>
		<key>com.apple.carkit.certificationOverrides.service</key>
		<true/>
		<key>com.apple.carkit.dnd.service</key>
		<true/>
		<key>com.apple.carkit.headUnitPairingPrompt.service</key>
		<true/>
		<key>com.apple.carkit.monitoring.service</key>
		<true/>
		<key>com.apple.carkit.navigation.service</key>
		<true/>
		<key>com.apple.carkit.navowners.service</key>
		<true/>
		<key>com.apple.carkit.pairing.service</key>
		<true/>
		<key>com.apple.carkit.reconnectiontime.service</key>
		<true/>
		<key>com.apple.carkit.remote-iap.service</key>
		<true/>
		<key>com.apple.carkit.service</key>
		<true/>
		<key>com.apple.carkit.service.lockdown</key>
		<true/>
		<key>com.apple.carkit.setupPromptDirector.service</key>
		<true/>
		<key>com.apple.carkit.theme-asset-library</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.CarKit.notifications</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/carkitd</string>
	<key>RemoteServices</key>
	<dict>
		<key>com.apple.carkit.remote-iap.service</key>
		<dict>
			<key>RequireEntitlement</key>
			<string>AppleInternal</string>
			<key>ServiceProperties</key>
			<dict>
				<key>UsesRemoteXPC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
