# beta/System/Library/LaunchDaemons/com.apple.tvremoted.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.tvremoted</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>tvremote.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>tvremote.ContinuityKeyBoardEnabled</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ContinuityKeyBoard.enabled</string>
			</dict>
			<key>tvremote.TapToSetupUsed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.Sharing.prefsChanged</string>
			</dict>
			<key>tvremote.airplayTVAdded</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.televisionAccessoryAdded</string>
			</dict>
			<key>tvremote.appleTVAdded</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.AppleTVAccessoryAdded</string>
			</dict>
			<key>tvremote.connectionRequested</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TVRemoteCore.connectionRequested</string>
			</dict>
			<key>tvremote.dataMigrationFinished</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.datamigrator.migrationDidFinish</string>
			</dict>
			<key>tvremote.mdmProfileInstalled</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ManagedConfiguration.profileListChanged</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.tvremotecore.xpc</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/tvremoted</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
