# beta/System/Library/LaunchDaemons/com.apple.fitcore.session.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MallocStackLoggingDirectory</key>
		<string>/var/mobile/Library/Seymour</string>
		<key>NSRunningFromLaunchd</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.fitcore.session</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>SeymourSession</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.SeymourSession</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.seymour.session</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/SeymourServices.framework/fitcoresessiond</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
