# beta/System/Library/LaunchDaemons/com.apple.appleidsetupd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.appleidsetupd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.authkit.user-info-changed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.authkit.user-info-changed</string>
			</dict>
			<key>com.apple.cdp.manatee-availability</key>
			<dict>
				<key>Name</key>
				<string>com.apple.cdp.manatee-availability</string>
			</dict>
			<key>com.apple.cdp.manatee-availability-changed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.cdp.manatee-availability-changed</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.devicemanagementclient.devicePostureChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.devicemanagementclient.devicePostureChanged</string>
			</dict>
			<key>com.apple.rapport.CompanionLinkDeviceAdded</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.rapport.CompanionLinkDeviceAdded</string>
			</dict>
			<key>com.apple.rapport.CompanionLinkDeviceRemoved</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.rapport.CompanionLinkDeviceRemoved</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.appleidsetup.xpc</key>
		<true/>
		<key>com.apple.appleidsetupd.agemigration.xpc</key>
		<true/>
		<key>com.apple.appleidsetupd.repair.xpc</key>
		<true/>
		<key>com.apple.appleidsetupd.setup.xpc</key>
		<true/>
		<key>com.apple.aps.appleidsetup</key>
		<true/>
		<key>com.apple.private.alloy.remoteappleidsetup-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.appleidsetup.notifications</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/appleidsetupd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
