# beta/System/Library/LaunchDaemons/com.apple.voicebankingd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>com.apple.coredata.CloudKitDebug</key>
		<string>0</string>
		<key>com.apple.coredata.XPCDebug</key>
		<string>0</string>
	</dict>
	<key>Label</key>
	<string>com.apple.voicebankingd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.voicebanking.systemtask.ServerConfigCheck</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.voicebanking.systemtask.V1ToV2VoiceMigrationTask</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>livepseech-cloud-prefs-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.accessibility.livespeech</string>
			</dict>
			<key>livepseech-local-prefs-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.livespeech.localprefschanged</string>
			</dict>
			<key>sirittstraining-task-cancelled</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SiriTTSTrainingAgent.taskEvent.cancelled</string>
			</dict>
			<key>sirittstraining-task-done</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SiriTTSTrainingAgent.taskEvent.done</string>
			</dict>
			<key>sirittstraining-task-failed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SiriTTSTrainingAgent.taskEvent.failed</string>
			</dict>
			<key>sirittstraining-task-paused</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SiriTTSTrainingAgent.taskEvent.event</string>
			</dict>
			<key>sirittstraining-task-running</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SiriTTSTrainingAgent.taskEvent.running</string>
			</dict>
			<key>sirittstraining-task-submitted</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SiriTTSTrainingAgent.taskEvent.submitted</string>
			</dict>
			<key>sirittstraining-task-undefined</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SiriTTSTrainingAgent.taskEvent.undefined</string>
			</dict>
			<key>system-first-unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.voicebankingd</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.voicebanking.usernotifications</key>
		<true/>
		<key>com.apple.voicebanking.services</key>
		<true/>
		<key>com.apple.voicebanking.store</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/TextToSpeechVoiceBankingSupport.framework/Support/voicebankingd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
