# beta/System/Library/LaunchDaemons/com.apple.osanalytics.osanalyticshelper.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.osanalytics.osanalyticshelper</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.jobquiescence</key>
		<dict>
			<key>com.apple.osanalytics.jobquiescence</key>
			<dict>
				<key>SubscriberName</key>
				<string>osanalyticshelper</string>
			</dict>
		</dict>
		<key>com.apple.remoteservicediscovery.events</key>
		<dict>
			<key>com.apple.SubmitDiagInfo.available</key>
			<dict>
				<key>AvailableService</key>
				<string>com.apple.osanalytics.logRelay</string>
				<key>DeviceType</key>
				<string>ncm-host</string>
			</dict>
			<key>com.apple.SubmitDiagInfo.available.display</key>
			<dict>
				<key>AvailableService</key>
				<string>com.apple.osanalytics.logRelay</string>
				<key>DeviceType</key>
				<string>display-host</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.osanalytics.daily.telemetry</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>3600</integer>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.osanalytics.postboot.telemetry</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>300</integer>
				<key>GracePeriod</key>
				<integer>300</integer>
			</dict>
			<key>com.apple.osanalytics.stability-monitor</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>21600</integer>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>21600</integer>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.osanalytics.state-monitor</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>21600</integer>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>21600</integer>
				<key>PowerNap</key>
				<true/>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.osanalytics.submissions.cleanup</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>43200</integer>
				<key>Interval</key>
				<integer>43200</integer>
				<key>PowerNap</key>
				<true/>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.osanalytics.submissions.submit</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>3600</integer>
				<key>PowerNap</key>
				<true/>
				<key>Repeating</key>
				<true/>
				<key>RequireInexpensiveNetworkConnectivity</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.osanalytics.submissions.submit-permissive</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Interval</key>
				<integer>43200</integer>
				<key>PowerNap</key>
				<true/>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>LimitLoadToSessionType</key>
	<dict>
		<key>#IfFeatureFlagEnabled</key>
		<string>UserManagement/SystemSessionD1</string>
		<key>#Then</key>
		<string>System</string>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.osanalytics.diagnosticmonitor</key>
		<true/>
		<key>com.apple.osanalytics.osanalyticshelper</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/osanalyticshelper</string>
		<string>server-init</string>
	</array>
	<key>RemoteServices</key>
	<dict>
		<key>com.apple.osanalytics.logTransfer</key>
		<dict>
			<key>ExposedToUntrustedDevices</key>
			<true/>
			<key>LimitExposedFromDeviceType</key>
			<array>
				<string>compute-node</string>
			</array>
			<key>RequireEntitlement</key>
			<string>com.apple.ReportCrash.antenna-access</string>
			<key>ServiceProperties</key>
			<dict>
				<key>UsesRemoteXPC</key>
				<true/>
			</dict>
		</dict>
	</dict>
</dict>
</plist>
