# beta/System/Library/LaunchDaemons/com.apple.batteryintelligenced.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.batteryintelligenced</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.coreduetcontext.client_event_stream</key>
		<dict>
			<key>notification</key>
			<dict>
				<key>client-identifier</key>
				<string>com.apple.batteryintelligenced.contextstore-registration</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.powerui.mclstatuschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.powerui.mclstatuschanged</string>
			</dict>
			<key>com.apple.system.powersources.chargingiconography</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powersources.chargingiconography</string>
			</dict>
			<key>com.apple.system.powersources.source</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powersources.source</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.batteryintelligenced.dailytask</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.batteryintelligenced.batteryanalysis</key>
		<true/>
		<key>com.apple.batteryintelligenced.batteryanalysismanager</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.batteryintelligence-notifications</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Background</string>
	<key>Program</key>
	<string>/usr/libexec/batteryintelligenced</string>
	<key>UserName</key>
	<string>root</string>
</dict>
</plist>
