# beta/System/Library/LaunchDaemons/com.apple.contacts.donation-agent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.contacts.donation-agent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.contactsd.analysis</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<real>86400</real>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.contactsd.managed-duplicate-service</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<real>86400</real>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.contacts.donation.agent</key>
		<true/>
		<key>com.apple.contacts.donation.maintenance</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/ContactsDonation.framework/Versions/A/Support/contactsdonationagent</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
