# beta/System/Library/LaunchDaemons/com.apple.accessoryd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.accessoryd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.accessoryd.btServiceMatching</key>
		<dict>
			<key>com.apple.accessoryd.matching.btService</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.accessoryd.digitalIDMatching</key>
		<dict>
			<key>com.apple.accessoryd.matching.digitalID</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.accessoryd.nf-events</key>
		<dict>
			<key>com.apple.accessoryd.nf-events</key>
			<dict/>
		</dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.accessoryd.matching.A2868.billboard</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOParentMatch</key>
				<dict>
					<key>IOPropertyMatch</key>
					<dict>
						<key>idProduct</key>
						<integer>5904</integer>
						<key>idVendor</key>
						<integer>1452</integer>
					</dict>
					<key>IOProviderClass</key>
					<string>IOUSBHostInterface</string>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleUSBHostBillboardDevice</string>
			</dict>
			<key>com.apple.accessoryd.matching.A2869.billboard</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOParentMatch</key>
				<dict>
					<key>IOPropertyMatch</key>
					<dict>
						<key>idProduct</key>
						<integer>5905</integer>
						<key>idVendor</key>
						<integer>1452</integer>
					</dict>
					<key>IOProviderClass</key>
					<string>IOUSBHostInterface</string>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleUSBHostBillboardDevice</string>
			</dict>
			<key>com.apple.accessoryd.matching.OOBPairing</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IOAccessoryOOBPairingInterface</string>
			</dict>
			<key>com.apple.accessoryd.matching.authCP</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>AppleAuthCP</string>
			</dict>
			<key>com.apple.accessoryd.matching.usbHost</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>SupportsIAP</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>IOUSBHostInterface</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.nfcacd.multitag.state.change</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nfcacd.multitag.state.change</string>
			</dict>
		</dict>
		<key>com.apple.timed.accessory-timesync</key>
		<dict>
			<key>com.apple.timed.accessory-timesync</key>
			<dict/>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.accessories.assistivetouch</key>
		<true/>
		<key>com.apple.accessories.audio</key>
		<true/>
		<key>com.apple.accessories.blepairing</key>
		<true/>
		<key>com.apple.accessories.communications</key>
		<true/>
		<key>com.apple.accessories.connection-info-server</key>
		<true/>
		<key>com.apple.accessories.externalaccessory-server</key>
		<true/>
		<key>com.apple.accessories.hid</key>
		<true/>
		<key>com.apple.accessories.medialibrary</key>
		<true/>
		<key>com.apple.accessories.navigation</key>
		<true/>
		<key>com.apple.accessories.now-playing</key>
		<true/>
		<key>com.apple.accessories.oobbtpairing</key>
		<true/>
		<key>com.apple.accessories.transport-server</key>
		<true/>
		<key>com.apple.accessories.voiceover</key>
		<true/>
		<key>com.apple.accessoryd.nf-events</key>
		<true/>
		<key>com.apple.iap2d.ExternalAccessory.distributednotification.server</key>
		<dict>
			<key>HideUntilCheckIn</key>
			<true/>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.iap2d.distributednotification.server</key>
		<dict>
			<key>HideUntilCheckIn</key>
			<true/>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.iap2d.xpc</key>
		<true/>
		<key>com.apple.iaptransportd.xpc</key>
		<true/>
		<key>com.apple.timed.accessory-timesync</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CoreAccessories.framework/Support/accessoryd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
