# beta/System/Library/LaunchDaemons/com.apple.lsd.system.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<dict>
		<key>#Else</key>
		<false/>
		<key>#IfFeatureFlagDisabled</key>
		<string>LaunchServices/system_session_lsd</string>
		<key>#Then</key>
		<true/>
	</dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.lsd.system</string>
	<key>LimitLoadToSessionType</key>
	<string>System</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.lsd.system.mapdb</key>
		<true/>
		<key>com.apple.lsd.system.modifydb</key>
		<true/>
		<key>com.apple.lsd.system.rebuild</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/lsd</string>
		<string>systemSessionDaemon</string>
		<string>runAsRoot</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>1</integer>
</dict>
</plist>
