# beta/System/Library/LaunchDaemons/com.apple.textunderstandingd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.textunderstandingd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.suggestions.settingsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.suggestions.settingsChanged</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.TextUnderstanding.DocumentUnderstandingHarvesting</key>
		<true/>
		<key>com.apple.TextUnderstanding.DocumentUpdateHandler</key>
		<true/>
		<key>com.apple.TextUnderstanding.SmartReplies.Inference</key>
		<true/>
		<key>com.apple.TextUnderstanding.process</key>
		<true/>
		<key>com.apple.corespotlight.receiver.textunderstandingd</key>
		<true/>
		<key>com.apple.corespotlight.scheduled.receiver.textunderstandingd</key>
		<true/>
		<key>com.apple.textunderstandingd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/textunderstandingd</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>3</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
