# beta/System/Library/LaunchDaemons/com.apple.aned.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.aned</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationStateChanged</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationStateChanged</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.aned.exit</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>aned</key>
					<false/>
				</dict>
				<key>IOProviderClass</key>
				<string>IOResources</string>
			</dict>
			<key>com.apple.aned.launch</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>aned</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>IOResources</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.appleneuralengine</key>
		<true/>
		<key>com.apple.appleneuralengine.private</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/aned</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/aned</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
	<key>com.apple.ane.iokit-user-access</key>
	<true/>
</dict>
</plist>
