# beta/System/Library/LaunchDaemons/com.apple.icloud.fmflocatord.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.icloud.fmflocatord</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AppleDatePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleDatePreferencesChangedNotification</string>
			</dict>
			<key>AppleKeyboardsPreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleKeyboardsPreferencesChangedNotification</string>
			</dict>
			<key>AppleLanguagePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleLanguagePreferencesChangedNotification</string>
			</dict>
			<key>AppleNumberPreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleNumberPreferencesChangedNotification</string>
			</dict>
			<key>AppleTimePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleTimePreferencesChangedNotification</string>
			</dict>
			<key>com.apple.language.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.language.changed</string>
			</dict>
			<key>com.apple.locationd.authorization</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.authorization</string>
			</dict>
			<key>com.apple.locationd/Prefs</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd/Prefs</string>
			</dict>
			<key>com.apple.mobileme.fmf1.allowFindMyFriendsModification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobileme.fmf1.allowFindMyFriendsModification</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidunpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidunpair</string>
			</dict>
			<key>com.apple.purplebuddy.setupdone</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.purplebuddy.setupdone</string>
			</dict>
			<key>com.apple.system.hostname</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.hostname</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.icloud.fmflocatord.launch-on-reboot</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>120</integer>
				<key>GracePeriod</key>
				<integer>30</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.icloud.fmflocatord</key>
		<true/>
		<key>com.apple.icloud.fmflocatord.APSTokenWatcher.dev</key>
		<true/>
		<key>com.apple.icloud.fmflocatord.APSTokenWatcher.prod</key>
		<true/>
		<key>com.apple.icloud.fmflocatord.aps-demo</key>
		<true/>
		<key>com.apple.icloud.fmflocatord.aps-development</key>
		<true/>
		<key>com.apple.icloud.fmflocatord.aps-production</key>
		<true/>
		<key>com.apple.private.alloy.fmflocator-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/fmflocatord</string>
	<key>Umask</key>
	<integer>23</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
