# beta/System/Library/LaunchDaemons/com.apple.assistantd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.assistantd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>CKPersistentStorageKeysAndValuesUpdated.SiriUserCorrections</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysAndValuesUpdated.SiriUserCorrections</string>
			</dict>
			<key>CKPersistentStorageKeysAndValuesUpdated.synched</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysAndValuesUpdated.synched</string>
			</dict>
			<key>CKPersistentStorageKeysRemoved.SiriUserCorrections</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysRemoved.SiriUserCorrections</string>
			</dict>
			<key>CKPersistentStorageKeysRemoved.synched</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysRemoved.synched</string>
			</dict>
			<key>FMFDevicesChangedNotification</key>
			<dict>
				<key>Name</key>
				<string>FMFDevicesChangedNotification</string>
			</dict>
			<key>FMFMeDeviceChangedNotification</key>
			<dict>
				<key>Name</key>
				<string>FMFMeDeviceChangedNotification</string>
			</dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AppleDatePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleDatePreferencesChangedNotification</string>
			</dict>
			<key>MPStoreClientTokenDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>MPStoreClientTokenDidChangeNotification</string>
			</dict>
			<key>__ABDataBaseChangedByOtherProcessNotification</key>
			<dict>
				<key>Notification</key>
				<string>__ABDataBaseChangedByOtherProcessNotification</string>
			</dict>
			<key>com.apple.MobileAsset.AutoAssetNotification^com.apple.MobileAsset.MAAutoAsset^STARTUP_ACTIVATED</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.AutoAssetNotification^com.apple.MobileAsset.MAAutoAsset^STARTUP_ACTIVATED</string>
			</dict>
			<key>com.apple.MobileAsset.EmbeddedSpeech.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.EmbeddedSpeech.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.SpeechEndpointAssets.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.SpeechEndpointAssets.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.SpeechEndpointAssets.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.SpeechEndpointAssets.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceServices.VoiceResources.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceServices.VoiceResources.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssets.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssets.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssets.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssets.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssets.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssets.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssets.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssets.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsIPad.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsIPad.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsIPad.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsIPad.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsWatch.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsWatch.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerAssetsWatch.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerAssetsWatch.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerHSAssets.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerHSAssets.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerHSAssets.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerHSAssets.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerHSAssetsIPad.ma.cached-metadata-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerHSAssetsIPad.ma.cached-metadata-updated</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceTriggerHSAssetsIPad.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceTriggerHSAssetsIPad.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.SafariShared.Assistant.reload_plugin</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SafariShared.Assistant.reload_plugin</string>
			</dict>
			<key>com.apple.VideoSubscriberAccount.DidRegisterSubscription</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.VideoSubscriberAccount.DidRegisterSubscription</string>
			</dict>
			<key>com.apple.assistant.app_vocabulary</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.app_vocabulary</string>
			</dict>
			<key>com.apple.assistant.siri_settings_did_change</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.siri_settings_did_change</string>
			</dict>
			<key>com.apple.assistant.sync_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.sync_data_changed</string>
			</dict>
			<key>com.apple.assistant.sync_homekit_now</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.sync_homekit_now</string>
			</dict>
			<key>com.apple.assistant.sync_homekit_urgent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.sync_homekit_urgent</string>
			</dict>
			<key>com.apple.assistant.sync_needed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.sync_needed</string>
			</dict>
			<key>com.apple.chatkit.groups.siri_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.chatkit.groups.siri_data_changed</string>
			</dict>
			<key>com.apple.geoservices.siri_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.geoservices.siri_data_changed</string>
			</dict>
			<key>com.apple.homed.multi-user-status-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.multi-user-status-changed</string>
			</dict>
			<key>com.apple.homed.user-cloud-share.wake.com.apple.siri.data</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.user-cloud-share.wake.com.apple.siri.data</string>
			</dict>
			<key>com.apple.homed.user-cloud-share.wake.com.apple.siri.zonesharing</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.user-cloud-share.wake.com.apple.siri.zonesharing</string>
			</dict>
			<key>com.apple.icloud.fmip.siri_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.fmip.siri_data_changed</string>
			</dict>
			<key>com.apple.managedconfiguration.effectivesettingschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.effectivesettingschanged</string>
			</dict>
			<key>com.apple.media.entities.siri_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.media.entities.siri_data_changed</string>
			</dict>
			<key>com.apple.media.podcasts.siri_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.media.podcasts.siri_data_changed</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
			<key>com.apple.nanoregistry.watchdidbecomeactive</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.watchdidbecomeactive</string>
			</dict>
			<key>com.apple.parsecd.bag</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.parsecd.bag</string>
			</dict>
			<key>com.apple.proactive.PersonalizationPortrait.namedEntitiesInvalidated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.proactive.PersonalizationPortrait.namedEntitiesInvalidated</string>
			</dict>
			<key>com.apple.siri.client.state.DynamiteClientState.siri_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.client.state.DynamiteClientState.siri_data_changed</string>
			</dict>
			<key>com.apple.siri.inference.audio-app-signals-update</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.inference.audio-app-signals-update</string>
			</dict>
			<key>com.apple.siri.power.PowerContextPolicy.updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.power.PowerContextPolicy.updated</string>
			</dict>
			<key>com.apple.siri.preheat.quiet</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.preheat.quiet</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.siri.understanding</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.siri.understanding</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.siri.understanding.nl.overrides</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.siri.understanding.nl.overrides</string>
			</dict>
			<key>com.apple.siri.vocabulary.contacts_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.vocabulary.contacts_changed</string>
			</dict>
			<key>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_ASR_ASSISTANT</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_ASR_ASSISTANT</string>
			</dict>
			<key>com.apple.voiceservices.notification.voice-update</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voiceservices.notification.voice-update</string>
			</dict>
			<key>com.apple.voicetrigger.EarlyDetect</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicetrigger.EarlyDetect</string>
			</dict>
			<key>com.apple.voicetrigger.PHSProfileModified</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicetrigger.PHSProfileModified</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>MessageCenterWakeUpRequest</key>
			<dict>
				<key>serviceType</key>
				<string>SCD.MessageCenter.remoteIntelligence</string>
				<key>type</key>
				<string>server</string>
			</dict>
			<key>WakeUpRequest</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.siri.wakeup</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.assistantd.fmip</key>
		<true/>
		<key>com.apple.assistant.analytics</key>
		<true/>
		<key>com.apple.assistant.analytics-observation</key>
		<true/>
		<key>com.apple.assistant.announcement_state.service</key>
		<true/>
		<key>com.apple.assistant.audio-session-assertion</key>
		<true/>
		<key>com.apple.assistant.client</key>
		<true/>
		<key>com.apple.assistant.crossdevice.election.service</key>
		<true/>
		<key>com.apple.assistant.device-context</key>
		<true/>
		<key>com.apple.assistant.dictation</key>
		<true/>
		<key>com.apple.assistant.flow_service</key>
		<true/>
		<key>com.apple.assistant.notification_service</key>
		<true/>
		<key>com.apple.assistant.request-dispatcher.uibridge-service</key>
		<true/>
		<key>com.apple.assistant.security</key>
		<true/>
		<key>com.apple.assistant.settings</key>
		<true/>
		<key>com.apple.assistant.sync</key>
		<true/>
		<key>com.apple.assistant.uibridge-service</key>
		<true/>
		<key>com.apple.assistantd.aps</key>
		<true/>
		<key>com.apple.assistantd.managedstorage</key>
		<true/>
		<key>com.apple.private.alloy.siri.device-idswake</key>
		<true/>
		<key>com.apple.private.alloy.siri.icloud-idswake</key>
		<true/>
		<key>com.apple.private.alloy.siri.location-idswake</key>
		<true/>
		<key>com.apple.private.alloy.siri.proxy-idswake</key>
		<true/>
		<key>com.apple.siri.analytics.assistant</key>
		<true/>
		<key>com.apple.siri.client_lite</key>
		<true/>
		<key>com.apple.siri.device_resolution</key>
		<true/>
		<key>com.apple.siri.external_request</key>
		<true/>
		<key>com.apple.siri.intelligence-tool-helper-service</key>
		<true/>
		<key>com.apple.siri.location</key>
		<true/>
		<key>com.apple.siri.orchestration.capabilities</key>
		<true/>
		<key>com.apple.siri.orchestration.intelligencesession</key>
		<true/>
		<key>com.apple.siri.orchestration.sirix</key>
		<true/>
		<key>com.apple.siri.pairedsync</key>
		<true/>
		<key>com.apple.siri.pommes_search_xpc_service</key>
		<true/>
		<key>com.apple.siri.tether</key>
		<true/>
		<key>com.apple.siri.vocabularyupdates</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.crossdevicearbitration.feedback.notification</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.siri.NotificationSource</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AssistantServices.framework/assistantd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
