# beta/System/Library/LaunchDaemons/com.apple.backupd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamMemoryLimit</key>
		<integer>8192</integer>
		<key>JetsamPriority</key>
		<integer>25</integer>
	</dict>
	<key>Label</key>
	<string>com.apple.backupd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.alarm</key>
		<true/>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AppleLanguagePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleLanguagePreferencesChangedNotification</string>
			</dict>
			<key>com.apple.security.octagon.trust-status-change</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.security.octagon.trust-status-change</string>
			</dict>
			<key>com.apple.security.secureobjectsync.circlechanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.security.secureobjectsync.circlechanged</string>
			</dict>
			<key>com.apple.security.secureobjectsync.viewschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.security.secureobjectsync.viewschanged</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.backupd.cellular</key>
			<dict>
				<key>Delay</key>
				<integer>3153600000</integer>
				<key>Repeating</key>
				<false/>
			</dict>
			<key>com.apple.backupd.onBattery</key>
			<dict>
				<key>Delay</key>
				<integer>3153600000</integer>
				<key>Repeating</key>
				<false/>
			</dict>
			<key>com.apple.backupd.wifi</key>
			<dict>
				<key>Delay</key>
				<integer>3153600000</integer>
				<key>Repeating</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.backupd</key>
		<true/>
		<key>com.apple.backupd.installcoordination</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/MobileBackup.framework/backupd</string>
	</array>
</dict>
</plist>
