# beta/System/Library/LaunchDaemons/com.apple.appprotectiond.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.appprotectiond</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.accounts.changes</key>
		<dict>
			<key>AppleAccountAddedOrChanged</key>
			<dict>
				<key>com.apple.account.AppleAccount</key>
				<dict>
					<key>ChangeTypes</key>
					<array>
						<integer>1</integer>
						<integer>2</integer>
						<integer>3</integer>
					</array>
				</dict>
			</dict>
			<key>IDMSAccountChanged</key>
			<dict>
				<key>com.apple.account.idms</key>
				<dict>
					<key>AllChanges</key>
					<true/>
				</dict>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>ApplicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>ApplicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AppleLanguagePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleLanguagePreferencesChangedNotification</string>
			</dict>
			<key>com.apple.ManagedConfiguration.managedAppsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ManagedConfiguration.managedAppsChanged</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.appprotectiond.extensioninfo</key>
		<true/>
		<key>com.apple.appprotectiond.extensionmonitor</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.appprotectiond.guard</key>
		<true/>
		<key>com.apple.appprotectiond.maintenance</key>
		<true/>
		<key>com.apple.appprotectiond.read</key>
		<true/>
		<key>com.apple.appprotectiond.viewsubjectinfo</key>
		<true/>
		<key>com.apple.appprotectiond.viewsubjectmonitor</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.appprotectiond.write</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/appprotectiond</string>
	</array>
	<key>PublishesEvents</key>
	<array>
		<string>com.apple.appprotection.systemappevents</string>
		<string>com.apple.appprotection.preferences.change.event</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
