# beta/System/Library/LaunchDaemons/com.apple.deviceinterfaced.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.deviceinterfaced</string>
	<key>LaunchEvents</key>
	<dict/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.deviceinterfaced.device-interface-client</key>
		<true/>
		<key>com.apple.deviceinterfaced.tadfu-transport-server</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/DeviceInterface.framework/Support/deviceinterfaced</string>
	<key>ProgramArguments</key>
	<array/>
</dict>
</plist>
