# beta/System/Library/LaunchDaemons/com.apple.usb.networking.addNetworkInterface.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.usb.networking.addNetworkInterface</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>network_interface</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>IONetworkInterface</string>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/addNetworkInterface</string>
		<string>-x</string>
		<string>-s</string>
		<string>-c</string>
		<string>AppleUSBNCMData</string>
		<string>-c</string>
		<string>AppleUSBEthernetController</string>
		<string>-c</string>
		<string>AppleUSBEthernetHost</string>
		<string>-c</string>
		<string>AppleUSBDeviceNCMData</string>
		<string>-c</string>
		<string>AppleUSBECMData</string>
		<string>-c</string>
		<string>AppleEthernetAquantiaAqtion</string>
		<string>-c</string>
		<string>BCM5701Enet</string>
		<string>-c</string>
		<string>AppleIntelI210</string>
		<string>-c</string>
		<string>IOUserNetworkEthernet</string>
	</array>
</dict>
</plist>
