# beta/System/Library/LaunchDaemons/com.apple.nanoregistryd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/var/mobile/Library/DeviceRegistry.state/KeepAlive</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.nanoregistryd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.MobileGestalt.MobileGestaltEvents</key>
		<dict>
			<key>com.apple.nanoregistry.mgDeviceNameChanged</key>
			<dict>
				<key>Queries</key>
				<array>
					<string>UserAssignedDeviceName</string>
				</array>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.nano.nanoregistry.paireddeviceregistry</key>
		<true/>
		<key>com.apple.nanoregistry.termsacknowledgementregistry</key>
		<true/>
		<key>com.apple.nanoregistry.versionrequest</key>
		<true/>
		<key>com.apple.private.alloy.bluetoothregistry-idswake</key>
		<true/>
		<key>com.apple.private.alloy.bluetoothregistry.cloud-idswake</key>
		<true/>
		<key>com.apple.private.alloy.bluetoothregistryclassa-idswake</key>
		<true/>
		<key>com.apple.private.alloy.bluetoothregistryclassc-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.usernotifications.nanoregistryd</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/nanoregistryd</string>
	</array>
	<key>TransactionTimeLimitEnabled</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
	<key>WorkingDirectory</key>
	<string>/var/mobile/Library/Caches</string>
</dict>
</plist>
