# beta/System/Library/LaunchDaemons/com.apple.CommCenter-noncellular.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>ExitTimeOut</key>
	<integer>14</integer>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.CommCenter</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.CommCenter</key>
		<true/>
		<key>com.apple.basebandd.xpc</key>
		<true/>
		<key>com.apple.commcenter.ari.rt.xpc</key>
		<true/>
		<key>com.apple.commcenter.atcs.xpc</key>
		<true/>
		<key>com.apple.commcenter.coretelephony.xpc</key>
		<true/>
		<key>com.apple.commcenter.cupolicy.xpc</key>
		<true/>
		<key>com.apple.commcenter.mobile-helper-cbupdateservice</key>
		<true/>
		<key>com.apple.commcenter.oswatchdog</key>
		<true/>
		<key>com.apple.commcenter.tmk.xpc</key>
		<true/>
		<key>com.apple.commcenter.xpc</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.SIMSetupUIService</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.SOSBuddy</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.telephony.BasebandServicesNotification</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/CoreTelephony.framework/Support/CommCenter</string>
	</array>
	<key>UserName</key>
	<string>_wireless</string>
	<key>_LimitLoadFromDeviceTree</key>
	<string>baseband</string>
</dict>
</plist>
