# beta/System/Library/LaunchDaemons/com.apple.alwaysonexclavesd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>AfterInitialDemand</key>
		<true/>
		<key>Crashed</key>
		<true/>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.apple.alwaysonexclavesd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.alwaysonexclaves</key>
		<true/>
	</dict>
	<key>Program</key>
	<string>/usr/libexec/alwaysonexclavesd</string>
	<key>_Conclave</key>
	<string>com.apple.conclave.AlwaysOnExclaves</string>
</dict>
</plist>
