# beta/System/Library/LaunchDaemons/com.apple.nearbyd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MallocStackLoggingDirectory</key>
		<string>/var/mobile/tmp/com.apple.nearbyd</string>
	</dict>
	<key>GroupName</key>
	<string>mobile</string>
	<key>KeepAlive</key>
	<dict>
		<key>Crashed</key>
		<true/>
	</dict>
	<key>Label</key>
	<string>com.apple.nearbyd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bluetooth.discovery</key>
		<dict>
			<key>PrecisionFindingDiscovery</key>
			<dict>
				<key>bleRSSIThresholdHint</key>
				<integer>-120</integer>
				<key>bleScanRate</key>
				<string>Background</string>
				<key>discoveryFlags</key>
				<array/>
				<key>discoveryTypes</key>
				<array/>
				<key>label</key>
				<string>PrecisionFindingDiscovery</string>
			</dict>
		</dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.nearbyd</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IONameMatch</key>
				<string>rose</string>
			</dict>
			<key>com.apple.nearbyd.iokit</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>RequireClient</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>AppleSPURoseDriver</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.bluetooth.daemonStarted</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.bluetooth.daemonStarted</string>
			</dict>
			<key>com.apple.bluetooth.state</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.bluetooth.state</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidunpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidunpair</string>
			</dict>
			<key>com.apple.system.lowpowermode</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.lowpowermode</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.nearbyd.xpc.BTRanging</key>
		<true/>
		<key>com.apple.nearbyd.xpc.beacon</key>
		<true/>
		<key>com.apple.nearbyd.xpc.companion</key>
		<true/>
		<key>com.apple.nearbyd.xpc.diagnostics</key>
		<true/>
		<key>com.apple.nearbyd.xpc.nearbyinteraction</key>
		<true/>
		<key>com.apple.nearbyd.xpc.nearbyinteraction.observer</key>
		<true/>
		<key>com.apple.nearbyd.xpc.ranging</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/nearbyd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
