# beta/System/Library/LaunchDaemons/com.apple.translationd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.translationd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.siri.uaf.com.apple.sequoia.asset</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.sequoia.asset</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.speech.automaticspeechrecognition</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.speech.automaticspeechrecognition</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.translation.text</key>
		<true/>
		<key>com.apple.translationd</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/Translation.framework/translationd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
