# beta/System/Library/LaunchDaemons/com.apple.arkitd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.apple.arkitd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.arkit.service.control</key>
		<true/>
		<key>com.apple.arkit.service.geoTracking</key>
		<true/>
		<key>com.apple.arkit.service.location</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/arkitd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
