# beta/System/Library/LaunchDaemons/com.apple.dockaccessoryd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.dockaccessoryd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bluetooth.discovery</key>
		<dict>
			<key>DockKitAccessory</key>
			<dict>
				<key>discoveryFlags</key>
				<array/>
				<key>discoveryTypes</key>
				<array>
					<string>DockKitAccessory</string>
				</array>
				<key>label</key>
				<string>DockKitAccessory</string>
			</dict>
		</dict>
		<key>com.apple.dockaccessoryd.pairsetup</key>
		<dict>
			<key>url</key>
			<dict/>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.dockaccessoryd</key>
		<true/>
		<key>com.apple.dockaccessoryd.cameracapture</key>
		<true/>
		<key>com.apple.dockaccessoryd.certification</key>
		<true/>
		<key>com.apple.dockaccessoryd.debug</key>
		<true/>
		<key>com.apple.dockaccessoryd.pairing</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/dockaccessoryd</string>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
