# diff/tvOS/System/Library/LaunchDaemons/tmp0wci6flf.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.cloudphotod</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.cloudphotod.backup</key>
			<dict>
				<key>InvolvedProcesses</key>
				<array>
					<string>cloudphotod</string>
					<string>nsurlsessiond</string>
					<string>cloudd</string>
					<string>com.apple.photos.VideoConversionService</string>
					<string>com.apple.photos.ImageConversionService</string>
					<string>assetsd</string>
				</array>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>300</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>UserRequestedBackupTask</key>
				<true/>
			</dict>
			<key>com.apple.cloudphotod.configurationrefresh</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>43200</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>34560</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.cloudphotod.initialdownload.fastpass</key>
			<dict>
				<key>FeatureCodes</key>
				<array>
					<integer>100</integer>
					<integer>200</integer>
				</array>
				<key>InvolvedProcesses</key>
				<array>
					<string>cloudphotod</string>
					<string>nsurlsessiond</string>
					<string>cloudd</string>
					<string>com.apple.photos.VideoConversionService</string>
					<string>com.apple.photos.ImageConversionService</string>
					<string>assetsd</string>
				</array>
				<key>NonRepeatingTask</key>
				<dict>
					<key>FastPassTask</key>
					<dict>
						<key>ProcessingTaskIdentifiers</key>
						<array>
							<string>com.apple.cloudphotod.sync.discretionary</string>
						</array>
					</dict>
				</dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>ProducedResultIdentifiers</key>
				<array>
					<string>com.apple.cloudphotod.initialdownload.assets</string>
					<string>com.apple.cloudphotod.initialdownload.finished</string>
					<string>com.apple.cloudphotod.initialdownload.fastpass</string>
				</array>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.cloudphotod.maintenance</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>43200</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.cloudphotod.periodic.upload.computestates</key>
			<dict>
				<key>InvolvedProcesses</key>
				<dict>
					<key>0</key>
					<string>assetsd</string>
				</dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>18000</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>14400</integer>
				</dict>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.photos.cloud.datarepair</key>
		<true/>
		<key>com.apple.cloudphotod</key>
		<true/>
		<key>com.apple.cloudphotod.daemon</key>
		<true/>
		<key>com.apple.photos.cloud.pushnotification</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/CloudPhotoLibrary.framework/Support/cloudphotod</string>
	<key>ProgramArguments</key>
	<array>
		<string>cloudphotod</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
