# diff/tvOS/System/Library/LaunchDaemons/tmp0tj5ds8_.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/var/mobile/Library/UserConfigurationProfiles/MCOutstandingNagComputation</key>
			<true/>
			<key>/var/mobile/Library/UserConfigurationProfiles/profiledRunning</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.managedconfiguration.profiled</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.ManagedSettings.effective-settings.changed</key>
		<dict>
			<key>com.apple.managedconfiguration.ManagedSettings-effective-settings-changed</key>
			<dict>
				<key>SettingsGroups</key>
				<array>
					<string>account</string>
					<string>airDrop</string>
					<string>allowedClient</string>
					<string>application</string>
					<string>appStore</string>
					<string>calculator</string>
					<string>camera</string>
					<string>carPlay</string>
					<string>cellular</string>
					<string>dateAndTime</string>
					<string>faceTime</string>
					<string>findMy</string>
					<string>gameCenter</string>
					<string>intelligence</string>
					<string>keyboard</string>
					<string>media</string>
					<string>messages</string>
					<string>news</string>
					<string>notification</string>
					<string>passcode</string>
					<string>privacy</string>
					<string>safari</string>
					<string>siri</string>
					<string>webContent</string>
				</array>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>FirstUnlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>GracePeriodChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.BiometricKit.passcodeGracePeriodChanged</string>
			</dict>
			<key>ProvisioningProfileInstalled</key>
			<dict>
				<key>Notification</key>
				<string>MISProvisioningProfileInstalled</string>
			</dict>
			<key>ProvisioningProfileRemoved</key>
			<dict>
				<key>Notification</key>
				<string>MISProvisioningProfileRemoved</string>
			</dict>
			<key>com.apple.LoginKit.isLoggedIn</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LoginKit.isLoggedIn</string>
			</dict>
			<key>com.apple.keystore.memento.effaced</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.keystore.memento.effaced</string>
			</dict>
			<key>com.apple.system.lowpowermode</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.lowpowermode</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.managedconfiguration.profiled</key>
		<true/>
		<key>com.apple.managedconfiguration.profiled.public</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/profiled</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>0</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
