# diff/tvOS/System/Library/LaunchDaemons/tmp1yutze62.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.contactsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CNContactStoreLimitedAccessDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>CNContactStoreLimitedAccessDidChangeNotification</string>
			</dict>
			<key>__ABDataBaseChangedByOtherProcessNotification</key>
			<dict>
				<key>Notification</key>
				<string>__ABDataBaseChangedByOtherProcessNotification</string>
			</dict>
			<key>com.apple.LaunchServices.database</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LaunchServices.database</string>
			</dict>
			<key>com.apple.TelephonyUtilities.DefaultAppRelayTelephonySettingChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TelephonyUtilities.DefaultAppRelayTelephonySettingChanged</string>
			</dict>
			<key>com.apple.contacts.clientDidDisplayFavorites</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.contacts.clientDidDisplayFavorites</string>
			</dict>
			<key>com.apple.datamigrator.migrationDidFinish</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.datamigrator.migrationDidFinish</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.contactsd.BackupConfigurationService</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.contactsd.PurgeChangeHistoryService</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
			<key>com.apple.contactsd.VerifyCoreSpotlightService</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.askto.client.com.apple.contacts.commlimits</key>
		<true/>
		<key>com.apple.contactsd</key>
		<true/>
		<key>com.apple.contactsd.contact-provider</key>
		<true/>
		<key>com.apple.contactsd.launch-services-proxy</key>
		<true/>
		<key>com.apple.contactsd.support</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/Contacts.framework/Support/contactsd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
