# diff/tvOS/System/Library/LaunchDaemons/tmp1ubkv_ws.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.hangreporter</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.fsevents.matching</key>
		<dict>
			<key>com.apple.hangreporter.spool.modified</key>
			<dict>
				<key>Path</key>
				<array>
					<string>/var/root/Library/Caches/hangtracerd/spool</string>
					<string>/var/mobile/Library/Caches/com.apple.backboardd/spool</string>
				</array>
			</dict>
		</dict>
	</dict>
	<key>LowPriorityIO</key>
	<true/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.hangreporter</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/hangreporter</string>
	</array>
</dict>
</plist>
