# diff/tvOS/System/Library/LaunchDaemons/tmp19zd1e5m.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.cfnetwork.cfnetworkagent</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.cfnetwork.cfnetworkagent</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/CFNetwork.framework/CFNetworkAgent</string>
	</array>
	<key>StartInterval</key>
	<integer>86500</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
