# diff/tvOS/System/Library/LaunchDaemons/tmp3_2t9geg.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.online-auth-agent.xpc</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>Application Installed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.mis.opportunistic-validation.boot</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Delay</key>
				<integer>120</integer>
				<key>GracePeriod</key>
				<integer>600</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.mis.profile-garbage-collection</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>1800</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.online-auth-agent.check-indeterminates</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
			<key>com.apple.online-auth-agent.denylist-update</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.online-auth-agent.xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/online-auth-agent</string>
</dict>
</plist>
