# diff/tvOS/System/Library/LaunchDaemons/tmp44vg5alm.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>GroupName</key>
	<string>_networkd</string>
	<key>Label</key>
	<string>com.apple.symptomsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.network_change</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.config.network_change</string>
			</dict>
			<key>com.apple.springboard.SBApplicationNotificationStateChanged</key>
			<dict>
				<key>Notification</key>
				<string>SBApplicationNotificationStateChanged</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.symptomsd.task.cleanup-ndfdevicerecords</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Delay</key>
				<integer>86400</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequireScreenSleep</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.resource_notify.host</key>
		<dict>
			<key>HostSpecialPort</key>
			<integer>27</integer>
		</dict>
		<key>com.apple.symptom_analytics</key>
		<true/>
		<key>com.apple.symptoms.symptomsd.managed_events</key>
		<true/>
		<key>com.apple.usymptomsd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/symptomsd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/symptomsd</string>
	</array>
	<key>UserName</key>
	<string>_networkd</string>
	<key>WatchPaths</key>
	<array>
		<string>/Library/Preferences/com.apple.symptoms.plist</string>
	</array>
	<key>_LimitLoadFromVariant</key>
	<string>IsDarwinOS</string>
</dict>
</plist>
