# diff/tvOS/System/Library/LaunchDaemons/tmp5o6ej99p.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>Label</key>
	<string>com.apple.audiomxd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.BTAudioHALPluginAccessories</key>
		<true/>
		<key>com.apple.airplay.agent.services</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.airplay.discoverybroker</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.airplay.receiver.services</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.audio.AURemoteIOServer</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.audio.AudioComponentPrefs</key>
		<true/>
		<key>com.apple.audio.AudioComponentRegistrar</key>
		<true/>
		<key>com.apple.audio.AudioFileServer</key>
		<true/>
		<key>com.apple.audio.AudioQueueServer</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.audio.AudioSession</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.audio.AudioUnitServer</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.audio.PhaseXPCServer</key>
		<true/>
		<key>com.apple.audio.SystemSoundServer-iOS</key>
		<true/>
		<key>com.apple.audio.adam.xpc</key>
		<true/>
		<key>com.apple.audio.driver-registrar</key>
		<true/>
		<key>com.apple.audio.hapticd</key>
		<true/>
		<key>com.apple.audio.isolated.client.service</key>
		<true/>
		<key>com.apple.audio.voicetrigger.xpc</key>
		<true/>
		<key>com.apple.audioanchorservice</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>CoreAudioServices/SpatialAudioExperience</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.audiomxd.watchdog</key>
		<true/>
		<key>com.apple.carkit.sessionRequestHandler</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coreaudio.adam.hae.notification</key>
		<true/>
		<key>com.apple.coremedia.audiomxd.audioprocessingtap.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.bufferedairplayglobalroutingregistry.xpc</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.carplayavvc.xpc</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.endpoint.xpc</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.endpointplaybacksession.xpc</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.endpointrcgroup.xpc</key>
		<true/>
		<key>com.apple.coremedia.endpointremotecontrolsession.xpc</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.endpointstream.xpc</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.endpointstreamaudioengine.xpc</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.endpointuiagent.xpc</key>
		<true/>
		<key>com.apple.coremedia.nerotransportconnectionxpc</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.routediscoverer.xpc</key>
		<true/>
		<key>com.apple.coremedia.routingcontext.xpc</key>
		<true/>
		<key>com.apple.coremedia.routingsessionmanager.xpc</key>
		<true/>
		<key>com.apple.coremedia.samplebufferconsumer.xpc</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.sts</key>
		<true/>
		<key>com.apple.coremedia.systemcontroller.xpc</key>
		<true/>
		<key>com.apple.coremedia.virtualdisplaysession</key>
		<true/>
		<key>com.apple.coremedia.virtualframebuffer.xpc</key>
		<true/>
		<key>com.apple.coremedia.volumecontroller.xpc</key>
		<true/>
		<key>com.apple.mediaexperience.carplaymodecontroller.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.mediaexperience.endpoint.xpc</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>AirPlay/airplayd</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.mediaexperience.xctest.xpc</key>
		<true/>
		<key>com.apple.mediasafetynet.exceptions</key>
		<true/>
		<key>com.apple.mediasafetynet.pill</key>
		<true/>
		<key>com.apple.polarisaudio.workloadcontroller.xpc</key>
		<true/>
		<key>com.apple.roomcongruenceservice</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>CoreAudioServices/SpatialAudioExperience</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.soundanalysisd</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>SoundAnalysis/standaloneDaemonDisabled</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.streamcontrolservice</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>CoreAudioServices/SpatialAudioExperience</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.usernotifications.delegate.com.apple.coreaudio.adam.hae</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.coreaudio.adam.hae-watchos</key>
		<true/>
		<key>com.apple.virtualaudio.cider</key>
		<true/>
		<key>com.apple.voicetrigger.voicetriggerservice</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/audiomxd</string>
	</array>
	<key>PublishesEvents</key>
	<array>
		<string>com.apple.audio.piper</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>UserName</key>
	<string>mobile</string>
	<key>_Conclave</key>
	<string>com.apple.audiomxd.conclave</string>
</dict>
</plist>
