# diff/tvOS/System/Library/LaunchDaemons/tmp55qnejsn.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.usernotificationsd</string>
	<key>LaunchEvents</key>
	<dict/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.bulletinboard.settingspersistenceconnection</key>
		<true/>
		<key>com.apple.private.alloy.bulletinboard-idswake</key>
		<true/>
		<key>com.apple.private.alloy.continuity.notifications-idswake</key>
		<true/>
		<key>com.apple.usernotifications.alert-coordination</key>
		<true/>
		<key>com.apple.usernotifications.coreservice</key>
		<true/>
		<key>com.apple.usernotifications.listener</key>
		<true/>
		<key>com.apple.usernotifications.settings-vendor</key>
		<true/>
		<key>com.apple.usernotifications.tool</key>
		<true/>
		<key>com.apple.usernotificationsvendor.listener</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/UserNotificationsCore.framework/Support/usernotificationsd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
