# diff/tvOS/System/Library/LaunchDaemons/tmp6zgpk_va.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/private/var/run/com.apple.security.cryptexd/init</key>
			<false/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.security.cryptexd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.security.cryptex.xpc</key>
		<true/>
		<key>com.apple.security.cryptexd</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/usr/libexec/cryptexd</string>
	<key>PublishesEvents</key>
	<string>com.apple.security.cryptex.events</string>
	<key>RemoteServices</key>
	<dict>
		<key>com.apple.security.cryptexd.remote</key>
		<dict>
			<key>ExposedToUntrustedDevices</key>
			<array>
				<string>PreSetup</string>
				<string>HasInternalSecurityPolicies</string>
			</array>
			<key>LimitExposedFromDeviceType</key>
			<array>
				<string>compute-node</string>
			</array>
			<key>RequireEntitlement</key>
			<string>com.apple.private.security.cryptexd.remote</string>
			<key>ServiceProperties</key>
			<dict>
				<key>Features</key>
				<array>
					<string>CryptexInstall</string>
					<string>Cryptex1</string>
					<string>ReadIdentifiers</string>
					<string>Cryptex1,UseProductClass</string>
				</array>
				<key>ServiceVersion</key>
				<integer>3</integer>
				<key>UsesRemoteXPC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>_LimitLoadFromBootMode</key>
	<array>
		<string>rem</string>
	</array>
</dict>
</plist>
