# diff/tvOS/System/Library/LaunchDaemons/tmp6v61v2xb.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.hangtracerd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.hangtracerd.htse_state_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.hangtracerd.htse_state_changed</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.hangtracerd</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.hangtracerd.usernotifications</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/hangtracerd</string>
	</array>
</dict>
</plist>
