# diff/tvOS/System/Library/LaunchDaemons/tmp6ctx4vwh.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.perfdiagsselfenabled</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.perfdiagsselfenabled</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>ExpectedDuration</key>
				<integer>120</integer>
				<key>GracePeriod</key>
				<integer>10800</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>86400</integer>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>LowPriorityIO</key>
	<true/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.perfdiagnosticsselfenablement</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/perfdiagsselfenabled</string>
	</array>
</dict>
</plist>
