# diff/tvOS/System/Library/LaunchDaemons/tmp621_ch6p.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/var/mobile/Library/Preferences/com.apple.MobileStoreDemo.plist</key>
			<true/>
			<key>/var/root/Library/Preferences/com.apple.MobileStoreDemo.plist</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.mobilestoredemod</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.fsevents.matching</key>
		<dict>
			<key>DemoSettingsManagedPreferencesChanged</key>
			<dict>
				<key>Path</key>
				<string>/private/var/Managed Preferences/mobile/com.apple.demo-settings.plist</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.mobilestoredemoport</key>
		<true/>
		<key>com.apple.mobilestoredemod</key>
		<true/>
		<key>com.apple.msdtest.service</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/demod</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
