# diff/tvOS/System/Library/LaunchDaemons/tmp844dgquq.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.amsengagementd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.itunescloud.ICMusicSubscriptionStatusDidChangeDistributedNotification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.itunescloud.ICMusicSubscriptionStatusDidChangeDistributedNotification</string>
			</dict>
			<key>com.apple.nfcd.ams.accessory</key>
			<dict>
				<key>Name</key>
				<string>com.apple.nfcd.ams.accessory</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>SignificantTimeChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>SignificantTimeChangeNotification</string>
			</dict>
			<key>com.apple.appstored.ActivitySubEntitlementsCacheUpdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appstored.ActivitySubEntitlementsCacheUpdated</string>
			</dict>
			<key>com.apple.appstored.AppStoreSubEntitlementsCacheUpdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appstored.AppStoreSubEntitlementsCacheUpdated</string>
			</dict>
			<key>com.apple.appstored.HWBundleSubEntitlementsCacheUpdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appstored.HWBundleSubEntitlementsCacheUpdated</string>
			</dict>
			<key>com.apple.appstored.MusicSubEntitlementsCacheUpdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appstored.MusicSubEntitlementsCacheUpdated</string>
			</dict>
			<key>com.apple.appstored.NewsSubEntitlementsCacheUpdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appstored.NewsSubEntitlementsCacheUpdated</string>
			</dict>
			<key>com.apple.appstored.PodcastSubEntitlementsCacheUpdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appstored.PodcastSubEntitlementsCacheUpdated</string>
			</dict>
			<key>com.apple.appstored.TVSubEntitlementsCacheUpdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appstored.TVSubEntitlementsCacheUpdated</string>
			</dict>
			<key>com.apple.appstored.iCloudSubEntitlementsCacheUpdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appstored.iCloudSubEntitlementsCacheUpdated</string>
			</dict>
			<key>com.apple.family.family_updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.family.family_updated</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.amsengagementd.daily</key>
			<dict>
				<key>Interval</key>
				<real>86400</real>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.amsengagementd</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.AppleMediaServicesUI.engagement.notifications</key>
		<true/>
		<key>com.apple.xpc.amsengagementd</key>
		<true/>
		<key>com.apple.xpc.amsmediacatalogsync</key>
		<true/>
		<key>com.apple.xpc.amsserverdatacacheservice</key>
		<true/>
		<key>com.apple.xpc.amsuniversallinks</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AppleMediaServicesUI.framework/amsengagementd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
