# diff/tvOS/System/Library/LaunchDaemons/tmp8zgg5upk.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.imnotificationserviceagent</string>
	<key>LaunchEvents</key>
	<dict/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.private.alloy.messagenotification-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/IMDPersistence.framework/IMNotificationServiceAgent.app/IMNotificationServiceAgent</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
