# diff/tvOS/System/Library/LaunchDaemons/tmp_hkq3v63.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.accessoryd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.accessoryd.btServiceMatching</key>
		<dict>
			<key>com.apple.accessoryd.matching.btService</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.accessoryd.matching.usbHost</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>SupportsIAP</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>IOUSBHostInterface</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.accessories.assistivetouch</key>
		<true/>
		<key>com.apple.accessories.connection-info-server</key>
		<true/>
		<key>com.apple.accessories.externalaccessory-server</key>
		<true/>
		<key>com.apple.accessories.hid</key>
		<true/>
		<key>com.apple.accessories.medialibrary</key>
		<true/>
		<key>com.apple.accessories.now-playing</key>
		<true/>
		<key>com.apple.accessories.transport-server</key>
		<true/>
		<key>com.apple.accessories.voiceover</key>
		<true/>
		<key>com.apple.iap2d.xpc</key>
		<true/>
		<key>com.apple.iaptransportd.xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CoreAccessories.framework/Support/accessoryd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
