# diff/tvOS/System/Library/LaunchDaemons/tmp_iaehp7f.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.featureaccessd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.featureaccess.reload</key>
			<dict>
				<key>AllowBattery</key>
				<string>1</string>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<string>1</string>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.ind.cloudfeatures</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CloudSubscriptionFeatures.framework/featureaccessd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
