# diff/tvOS/System/Library/LaunchDaemons/tmp_9egk3f0.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MallocStackLoggingDirectory</key>
		<string>/private/var/tmp/com.apple.intelligentroutingd</string>
	</dict>
	<key>GroupName</key>
	<string>mobile</string>
	<key>Label</key>
	<string>com.apple.intelligentroutingd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.intelligentroutingd.xpc.media</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/intelligentroutingd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
