# diff/tvOS/System/Library/LaunchDaemons/tmpagatf6bt.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.siriknowledged</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.cascade.xpc_event.setChange</key>
		<dict>
			<key>com.apple.siri.inference.EntityMatcher</key>
			<dict/>
		</dict>
		<key>com.apple.coreduetcontext.client_event_stream</key>
		<dict>
			<key>notification</key>
			<dict>
				<key>client-identifier</key>
				<string>com.apple.siriknowledge</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AFLanguageCodeDidChangeDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>AFLanguageCodeDidChangeDarwinNotification</string>
			</dict>
			<key>INVoocabularyChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>INVoocabularyChangedNotification</string>
			</dict>
			<key>__ABDataBaseChangedByOtherProcessNotification</key>
			<dict>
				<key>Notification</key>
				<string>__ABDataBaseChangedByOtherProcessNotification</string>
			</dict>
			<key>com.apple.CascadeSets.DonateNow</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.CascadeSets.DonateNow</string>
			</dict>
			<key>com.apple.MobileAsset.AutoAssetAtomicNotification^ATOMIC_INSTANCE_DOWNLOADED</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.AutoAssetAtomicNotification^ATOMIC_INSTANCE_DOWNLOADED</string>
			</dict>
			<key>com.apple.MobileAsset.AutoAssetAtomicNotification^ATOMIC_INSTANCE_ELIMINATED</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.AutoAssetAtomicNotification^ATOMIC_INSTANCE_ELIMINATED</string>
			</dict>
			<key>com.apple.MobileAsset.AutoAssetAtomicNotification^ATOMIC_INSTANCE_NO_ENTRIES</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.AutoAssetAtomicNotification^ATOMIC_INSTANCE_NO_ENTRIES</string>
			</dict>
			<key>com.apple.assistant.siri_settings_did_change</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.assistant.siri_settings_did_change</string>
			</dict>
			<key>com.apple.gms.availability.notification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.gms.availability.notification</string>
			</dict>
			<key>com.apple.homekit.sync-data-cache-updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homekit.sync-data-cache-updated</string>
			</dict>
			<key>com.apple.icloud.fmip.siri_data_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.fmip.siri_data_changed</string>
			</dict>
			<key>com.apple.language.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.language.changed</string>
			</dict>
			<key>com.apple.proactive.PersonalizationPortrait.namedEntitiesInvalidated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.proactive.PersonalizationPortrait.namedEntitiesInvalidated</string>
			</dict>
			<key>com.apple.sysdiagnose.sysdiagnoseStarted</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sysdiagnose.sysdiagnoseStarted</string>
			</dict>
			<key>kAFPreferencesDidChangeDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>kAFPreferencesDidChangeDarwinNotification</string>
			</dict>
		</dict>
		<key>com.apple.uservault.events</key>
		<dict>
			<key>config</key>
			<dict>
				<key>bundleID</key>
				<string>com.apple.siriknowledged</string>
				<key>resources</key>
				<dict>
					<key>com.apple.siri.entity-matcher</key>
					<dict/>
				</dict>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.siri.xpc_activity.uaf-post-upgrade</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>GracePeriod</key>
				<integer>300</integer>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
			<key>com.apple.siri.xpc_activity.uaf-subscription-maintenance</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>GracePeriod</key>
				<integer>14400</integer>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresClassC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.siri.conversation_context_module</key>
		<true/>
		<key>com.apple.siri.inference.EntityMatcher.admin</key>
		<true/>
		<key>com.apple.siri.reference_resolution_module</key>
		<true/>
		<key>com.apple.siri.uaf.service</key>
		<true/>
		<key>com.apple.siriknowledged</key>
		<true/>
		<key>com.apple.siriknowledged.aps</key>
		<true/>
		<key>com.apple.siriknowledged.koa.donate</key>
		<true/>
		<key>com.apple.siriknowledged.vocabulary.admin</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/siriknowledged</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/siriknowledged</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
