# diff/tvOS/System/Library/LaunchDaemons/tmpc7sb9z7r.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.generativeexperiencesd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.GenerativeExperiences.PeriodicTasks.OutOfSpaceFollowUp</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>14400</integer>
				</dict>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.GenerativeFunctions.PeriodicTasks.AvailabilityUpdateTask.Boot</key>
			<dict>
				<key>NonRepeatingTask</key>
				<dict>
					<key>ScheduleAfter</key>
					<integer>0</integer>
					<key>TrySchedulingBefore</key>
					<integer>0</integer>
				</dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.GenerativeFunctions.PeriodicTasks.InstrumentationUpload</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.GenerativeFunctions.PeriodicTasks.ReportAvailabilityTask.Daily</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>43200</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
		</dict>
		<key>com.apple.biome.compute.publisher</key>
		<dict>
			<key>com.apple.generativeexperiences.failuretrackingtrigger</key>
			<dict>
				<key>client</key>
				<string>com.apple.GenerativeFunctions.generativeexperiencesd</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>GenerativeExperiences.FailureTracking</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.generativeexperiences.failuretrackingtrigger</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict/>
		<key>com.apple.notifyd.matching</key>
		<dict/>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.generativeexperiences.ExternalPartnerCredentialStorage</key>
		<true/>
		<key>com.apple.generativeexperiences.FailureTracking</key>
		<true/>
		<key>com.apple.generativeexperiences.availability.internal</key>
		<true/>
		<key>com.apple.generativeexperiences.availabilityService</key>
		<true/>
		<key>com.apple.generativeexperiences.classification</key>
		<true/>
		<key>com.apple.generativeexperiences.corefollowup</key>
		<true/>
		<key>com.apple.generativeexperiences.externaltextcomposition</key>
		<dict>
			<key>#Else</key>
			<false/>
			<key>#IfFeatureFlagDisabled</key>
			<string>TextComposer/TextComposerDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.generativeexperiences.generativeexperiencessession</key>
		<true/>
		<key>com.apple.generativeexperiences.retriever</key>
		<true/>
		<key>com.apple.generativeexperiences.summarization</key>
		<true/>
		<key>com.apple.generativeexperiences.sysdiagnose</key>
		<true/>
		<key>com.apple.generativeexperiences.textcomposition</key>
		<dict>
			<key>#Else</key>
			<false/>
			<key>#IfFeatureFlagDisabled</key>
			<string>TextComposer/TextComposerDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/GenerativeExperiencesRuntime.framework/generativeexperiencesd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
