# diff/tvOS/System/Library/LaunchDaemons/tmpdvvhy909.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.bluetoothuserd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>CKIdentityUpdateNotification</key>
			<dict>
				<key>Notification</key>
				<string>CKIdentityUpdateNotification</string>
			</dict>
			<key>com.apple.bluetooth.daemonStarted</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.bluetooth.daemonStarted</string>
			</dict>
			<key>com.apple.icloudpairing.idslaunchnotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloudpairing.idslaunchnotification</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.bluetooth.cloud.settings</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.bluetooth.cloud.settings</string>
			</dict>
			<key>com.apple.system.hostname</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.hostname</string>
			</dict>
			<key>com.apple.xpc.activity</key>
			<dict>
				<key>com.apple.bluetoothuser.cloud.retryFetch</key>
				<dict>
					<key>AllowBattery</key>
					<true/>
					<key>Priority</key>
					<string>Utility</string>
					<key>Repeating</key>
					<false/>
					<key>RequireNetworkConnectivity</key>
					<true/>
					<key>RequiresClassC</key>
					<true/>
				</dict>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.bluetoothuserd</key>
		<true/>
		<key>com.apple.bluetoothuser.xpc</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.bluetoothuserd.UserNotification</key>
		<true/>
	</dict>
	<key>Program</key>
	<string>/usr/libexec/bluetoothuserd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/bluetoothuserd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>com.apple.alarm</key>
	<true/>
</dict>
</plist>
