# diff/tvOS/System/Library/LaunchDaemons/tmpd1gfidl0.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.assistivetouchd</string>
	<key>MachServices</key>
	<dict>
		<key>UIASTNotificationCenter</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.accessibility.AXAssistiveTouchServer</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.accessibility.AXSwitchControlServer</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.assistivetouchd</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.assistivetouchd.drag.xpc</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.private.alloy.accessibility.switchcontrol-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/AssistiveTouch.app/assistivetouchd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>ThrottleInterval</key>
	<integer>1</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
