# diff/tvOS/System/Library/LaunchDaemons/tmpdq4_q5qr.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.MobileFileIntegrity</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.fsevents.matching</key>
		<dict>
			<key>com.apple.MobileFileIntegrity.CodeRequirementPrefsChanged</key>
			<dict>
				<key>Path</key>
				<string>/Library/Preferences/com.apple.security.coderequirements.plist</string>
			</dict>
			<key>com.apple.MobileFileIntegrity.LibraryValidationPrefsChanged</key>
			<dict>
				<key>Path</key>
				<string>/Library/Preferences/com.apple.security.libraryvalidation.plist</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.security.amfid.cleanmetrics</key>
			<dict>
				<key>Delay</key>
				<integer>1209600</integer>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.security.amfid.reportmetrics</key>
			<dict>
				<key>Delay</key>
				<integer>86400</integer>
				<key>GracePeriod</key>
				<integer>3600</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.MobileFileIntegrity</key>
		<dict>
			<key>HostSpecialPort</key>
			<integer>18</integer>
		</dict>
		<key>com.apple.amfi.lockdown</key>
		<true/>
		<key>com.apple.amfi.nsxpc</key>
		<true/>
		<key>com.apple.amfi.xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/amfid</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
