# diff/tvOS/System/Library/LaunchDaemons/tmpe5eugwdr.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<dict>
		<key>#Else</key>
		<true/>
		<key>#IfFeatureFlagEnabled</key>
		<string>IntelligenceFlow/IntelligenceFlow</string>
		<key>#Then</key>
		<false/>
	</dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.intelligenceflowd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.biome.compute.publisher</key>
		<dict>
			<key>com.apple.intelligenceflowd.DeviceUpgrade</key>
			<dict>
				<key>client</key>
				<string>com.apple.intelligenceflow.intelligenceflowd</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>Device.Metadata</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.intelligenceflowd.DeviceUpgrade</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>TKToolkitDatabaseChangedNotification</key>
			<dict/>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.os-eligibility-domain.change.greymatter</key>
			<dict/>
			<key>com.apple.siri.uaf.tool_retrieval.denylist</key>
			<dict/>
			<key>com.apple.siri.uaf.tool_retrieval.embeddingModel</key>
			<dict/>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict/>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.intelligenceflow.internal</key>
		<true/>
		<key>com.apple.intelligenceflow.orchestrator</key>
		<true/>
		<key>com.apple.intelligenceflow.querydecoration</key>
		<true/>
		<key>com.apple.intelligenceflow.snippet-streaming</key>
		<true/>
		<key>com.apple.intelligenceflow.toolbox</key>
		<true/>
		<key>com.apple.intelligenceflow.transcript-entity-querying</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/IntelligenceFlowRuntime.framework/intelligenceflowd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
