# diff/tvOS/System/Library/LaunchDaemons/tmpej5e7qia.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.knowledgeconstructiond</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.FullPipeline</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>ResourceIntensive</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict/>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.intelligenceplatform.KnowledgeConstruction</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/IntelligencePlatformCore.framework/knowledgeconstructiond</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
