# diff/tvOS/System/Library/LaunchDaemons/tmpf0or5ej2.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MallocMaxMagazines</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.wifip2pd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.wifip2pd</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/wifip2pd</string>
</dict>
</plist>
