# diff/tvOS/System/Library/LaunchDaemons/tmpfsamov0w.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<dict>
		<key>#Else</key>
		<false/>
		<key>#IfFeatureFlagEnabled</key>
		<string>CoreMedia/MediaplaybackdPressuredExit</string>
		<key>#Then</key>
		<true/>
	</dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.mediaplaybackd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.coremedia.mediaplaybackd.asset.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.assetcacheinspector.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.assetdownloader.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.assetimagegenerator.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.audioprocessingtap.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.bytestream.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.controlcommands.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.coremediatestfoundation.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.cpe.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.cpeprotector.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.customurlloader.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.figcontentkeyboss.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.figcontentkeysession.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.figcpecryptor.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.figmetriceventtimeline.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.formatreader.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.mutablecomposition.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.mutablemovie.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.player.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.remaker.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.samplebufferaudiorenderer.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.samplebufferrendersynchronizer.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.samplegenerator.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.sandboxserver.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.videocompositor.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.videoqueue</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.coremedia.mediaplaybackd.videotarget.xpc</key>
		<true/>
		<key>com.apple.coremedia.mediaplaybackd.visualcontext.xpc</key>
		<true/>
		<key>com.apple.coremedia.neroid</key>
		<dict>
			<key>#IfFeatureFlagEnabled</key>
			<string>CoreMedia/SeparatePlayerProcess</string>
			<key>#Then</key>
			<true/>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/mediaplaybackd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
