# diff/tvOS/System/Library/LaunchDaemons/tmpflvgvjdh.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/private/var/containers/Shared/SystemGroup/systemgroup.com.apple.installcoordinationd/Library/InstallCoordination/PendingWork.plist</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.installcoordinationd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.InstallProgress.xpc</key>
		<true/>
		<key>com.apple.installcoordinationd</key>
		<true/>
		<key>com.apple.installcoordinationd.CacheDelete</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/InstallCoordination.framework/Support/installcoordinationd</string>
	</array>
	<key>PublishesEvents</key>
	<string>com.apple.InstallProgress.events</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
