# diff/tvOS/System/Library/LaunchDaemons/tmphjr5u1lf.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MallocNanoZone</key>
		<string>0</string>
	</dict>
	<key>Label</key>
	<string>com.apple.assistant_cdmd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.siri.VoiceShortcuts.DataDidUpdateNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.VoiceShortcuts.DataDidUpdateNotification</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.if.planner</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.if.planner</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.if.planner.overrides</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.if.planner.overrides</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.siri.understanding</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.siri.understanding</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.siri.understanding.nl.overrides</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.siri.understanding.nl.overrides</string>
			</dict>
			<key>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_NL</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_NL</string>
			</dict>
			<key>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_NL_OVERRIDES</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.SIRI_UNDERSTANDING_NL_OVERRIDES</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.siri.cdm.xpc_activity.maintenance</key>
			<dict>
				<key>Name</key>
				<string>com.apple.siri.cdm.xpc_activity.maintenance</string>
			</dict>
			<key>com.apple.siri.cdm.xpc_activity.post_install</key>
			<dict>
				<key>Name</key>
				<string>com.apple.siri.cdm.xpc_activity.post_install</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.assistant.cdm</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/ContinuousDialogManagerService.framework/assistant_cdmd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
