# diff/tvOS/System/Library/LaunchDaemons/tmphghec4ky.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AppIntents</key>
	<dict>
		<key>MetadataAbsolutePaths</key>
		<array>
			<string>/System/Library/PrivateFrameworks/IMDaemonCore.framework</string>
		</array>
	</dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>NSRunningFromLaunchd</key>
		<string>1</string>
	</dict>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamMemoryLimit</key>
		<integer>3000</integer>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.apple.imagent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>Application Installed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>Application Uninstalled</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
			<key>com.apple.authkit.user-info-changed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.authkit.user-info-changed</string>
			</dict>
			<key>com.apple.private.IMCore.LoggedIntoHSA2</key>
			<dict>
				<key>Name</key>
				<string>com.apple.private.IMCore.LoggedIntoHSA2</string>
			</dict>
			<key>com.apple.private.IMCore.LoggedIntoNonHSA2</key>
			<dict>
				<key>Name</key>
				<string>com.apple.private.IMCore.LoggedIntoNonHSA2</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict/>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.imagent</key>
		<true/>
		<key>com.apple.corespotlight.daemon.messages</key>
		<true/>
		<key>com.apple.hsa-authentication-server</key>
		<true/>
		<key>com.apple.imagent.cache-delete</key>
		<true/>
		<key>com.apple.imagent.embedded.auth</key>
		<true/>
		<key>com.apple.incoming-call-filter-server</key>
		<true/>
		<key>com.apple.madrid-idswake</key>
		<true/>
		<key>com.apple.madrid.lite-idswake</key>
		<true/>
		<key>com.apple.messages.critical-messaging</key>
		<true/>
		<key>com.apple.pairedsync.messagessync</key>
		<true/>
		<key>com.apple.private.appintents.delegate.com.apple.imagent</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.MobileSMS</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/IMCore.framework/imagent.app/imagent</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
	<key>WorkingDirectory</key>
	<string>/tmp</string>
	<key>_ExponentialThrottling</key>
	<true/>
	<key>_UrgentLogSubmission</key>
	<dict>
		<key>CheckedAllocations</key>
		<true/>
	</dict>
</dict>
</plist>
