# diff/tvOS/System/Library/LaunchDaemons/tmpi6o6thfc.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>CFBundleIdentifier</key>
	<string>com.apple.mediasetupd</string>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.mediasetupd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.mediasetupd.publicDBRefresh</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>300</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.homed.multi-user-status-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.multi-user-status-changed</string>
			</dict>
			<key>com.apple.homed.user-cloud-share.wake.com.apple.mediaservicesbroker.container</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.user-cloud-share.wake.com.apple.mediaservicesbroker.container</string>
			</dict>
			<key>com.apple.homehubd.endpointActivated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homehubd.endpointActivated</string>
			</dict>
			<key>com.apple.homehubd.endpointDeactivated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homehubd.endpointDeactivated</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.mediasetupd</key>
		<true/>
		<key>com.apple.mediasetupd.server</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/mediasetupd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
