# diff/tvOS/System/Library/LaunchDaemons/tmpicm9xc6y.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.apple.accessibility.axassetsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.LaunchServices.ApplicationsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LaunchServices.ApplicationsChanged</string>
			</dict>
			<key>com.apple.MobileAsset.TTSAXResourceModelAssets.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.TTSAXResourceModelAssets.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.UAF.Siri.TextToSpeech.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.UAF.Siri.TextToSpeech.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceServices.CustomVoice.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceServices.CustomVoice.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceServices.GryphonVoice.ma.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceServices.GryphonVoice.ma.new-asset-installed</string>
			</dict>
			<key>com.apple.SynthesisProvider.updatedVoices</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.SynthesisProvider.updatedVoices</string>
			</dict>
			<key>com.apple.coreaudio.components.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.coreaudio.components.changed</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.trial.NamespaceUpdate.SIRI_TEXT_TO_SPEECH</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.trial.NamespaceUpdate.SIRI_TEXT_TO_SPEECH</string>
			</dict>
			<key>com.apple.ttsasset.NewAssetNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ttsasset.NewAssetNotification</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.ax.ImageCaptionModel.LaunchActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.ax.MagnifierAsset.LaunchActivity</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>259200</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.accessibility.axassetsd.service</key>
		<true/>
		<key>com.apple.accessibility.voices</key>
		<true/>
		<key>com.apple.accessibility.voices.admin</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AXAssetLoader.framework/Support/axassetsd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
