# diff/tvOS/System/Library/LaunchDaemons/tmpjevg6ncf.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AbandonProcessGroup</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.memory-maintenance</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.memory-maintenance.rearm-exc-resource</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresUserInactivity</key>
				<false/>
			</dict>
			<key>com.apple.memory-maintenance.report-jetsam-telemetry</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>28800</integer>
				</dict>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.memory-maintenance.report-mach-memory-info</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<integer>86400</integer>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
			<key>com.apple.memory-maintenance.system-hwm</key>
			<dict>
				<key>MayRebootDevice</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Delay</key>
					<integer>0</integer>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>kBGSTProtectionClassD</string>
				<key>RequiresSignificantUserInactivity</key>
				<true/>
				<key>ResourceIntensive</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>LowPriorityIO</key>
	<true/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mmaintenanced</key>
		<true/>
		<key>com.apple.mmaintenanced.server</key>
		<dict>
			<key>HostSpecialPort</key>
			<integer>33</integer>
		</dict>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/mmaintenanced</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>com.apple.private.memoryinfo</key>
	<true/>
</dict>
</plist>
