# diff/tvOS/System/Library/LaunchDaemons/tmpl_ouo5bz.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/var/mobile/Library/homed/homeeventsd/reset</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.apple.homeeventsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.private.home.hindsight.dailyRecurringEventRecordUploadTask</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>ShouldWakeDevice</key>
				<true/>
			</dict>
			<key>com.apple.private.home.hindsight.onBootEventRecordUploadTask</key>
			<dict>
				<key>NonRepeatingTask</key>
				<dict>
					<key>TrySchedulingBefore</key>
					<integer>35</integer>
				</dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>ShouldWakeDevice</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.homekitevents.xpc</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/HomeKitEvents.framework/Support/homeeventsd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
