# diff/tvOS/System/Library/LaunchDaemons/tmpns5cbgkv.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>GroupName</key>
	<string>mobile</string>
	<key>Label</key>
	<string>com.apple.remoteappintentsd</string>
	<key>LaunchEvents</key>
	<dict/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.remoteappintentsd.appevent</key>
		<true/>
	</dict>
	<key>NSBonjourServices</key>
	<array>
		<string>_remoteAppIntents._udp</string>
		<string>_remoteAppIntents._tcp</string>
	</array>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/remoteappintentsd</string>
	</array>
	<key>Umask</key>
	<integer>54</integer>
	<key>UserName</key>
	<string>mobile</string>
	<key>com.apple.security.network.server</key>
	<true/>
</dict>
</plist>
