# diff/tvOS/System/Library/LaunchDaemons/tmpoo4b2f74.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.idleassetsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>ACDAccountStoreDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>ACDAccountStoreDidChangeNotification</string>
			</dict>
			<key>AppleLanguagePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleLanguagePreferencesChangedNotification</string>
			</dict>
			<key>com.apple.PineBoard.FinishedStartup</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.PineBoard.FinishedStartup</string>
			</dict>
			<key>com.apple.TVPhotos.preferences.notify</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TVPhotos.preferences.notify</string>
			</dict>
			<key>com.apple.TVScreenSaver.IdleScreenRefreshIntervalChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TVScreenSaver.IdleScreenRefreshIntervalChanged</string>
			</dict>
			<key>com.apple.TVScreenSaver.IdleScreenScreenSaverTypeChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TVScreenSaver.IdleScreenScreenSaverTypeChanged</string>
			</dict>
			<key>com.apple.TVScreenSaver.PhotosSharingFilterChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TVScreenSaver.PhotosSharingFilterChanged</string>
			</dict>
			<key>com.apple.TVScreenSaver.RequestMemoriesRefresh</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TVScreenSaver.RequestMemoriesRefresh</string>
			</dict>
			<key>com.apple.TVScreenSaver.TVScreenSaverAssetServiceManagerUpdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TVScreenSaver.TVScreenSaverAssetServiceManagerUpdated</string>
			</dict>
			<key>com.apple.TVScreenSaver.preferences.notify</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.TVScreenSaver.preferences.notify</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.idleassetsd.CacheDelete</key>
		<true/>
		<key>com.apple.idleassetsd.ambient.xpc</key>
		<true/>
		<key>com.apple.idleassetsd.xpc</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/idleassetsd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
