# diff/tvOS/System/Library/LaunchDaemons/tmpojg8zngk.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AppIntents</key>
	<dict>
		<key>MetadataAbsolutePaths</key>
		<array>
			<string>/System/Library/PrivateFrameworks/AppStoreComponentsDaemonKit.framework</string>
		</array>
	</dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.appstorecomponentsd</string>
	<key>LaunchEvents</key>
	<dict/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.appstorecomponentsd.xpc</key>
		<true/>
		<key>com.apple.chrono.event-service.appstorecomponentsd</key>
		<true/>
		<key>com.apple.private.appintents.delegate.com.apple.appstorecomponentsd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AppStoreComponents.framework/Support/appstorecomponentsd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
