# diff/tvOS/System/Library/LaunchDaemons/tmpp4e_s7hr.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.relatived</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.relatived.internal</key>
		<true/>
		<key>com.apple.relatived.public</key>
		<true/>
		<key>com.apple.relatived.status</key>
		<true/>
		<key>com.apple.relatived.tempest</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/relatived</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
