# diff/tvOS/System/Library/LaunchDaemons/tmppajrih16.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.amsaccountsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.nfcd.ams.accessory</key>
			<dict>
				<key>Name</key>
				<string>com.apple.nfcd.ams.accessory</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.AppleMediaServices.accountCachedData.dataUpdated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AppleMediaServices.accountCachedData.dataUpdated</string>
			</dict>
			<key>com.apple.AppleMediaServices.activeicloudaccountchanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AppleMediaServices.activeicloudaccountchanged</string>
			</dict>
			<key>com.apple.AppleMediaServices.eligibilityoverridechanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.AppleMediaServices.eligibilityoverridechanged</string>
			</dict>
			<key>com.apple.StoreServices.SSAccountStore.activeaccountchanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.StoreServices.SSAccountStore.activeaccountchanged</string>
			</dict>
			<key>com.apple.ams.privateListeningChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ams.privateListeningChanged</string>
			</dict>
			<key>com.apple.ams.provision-biometrics</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ams.provision-biometrics</string>
			</dict>
			<key>com.apple.homed.multi-user-status-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.multi-user-status-changed</string>
			</dict>
			<key>com.apple.homed.user-cloud-share.repair.wake.com.apple.applemediaservices.multiuser</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.user-cloud-share.repair.wake.com.apple.applemediaservices.multiuser</string>
			</dict>
			<key>com.apple.homed.user-cloud-share.wake.com.apple.applemediaservices.multiuser</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.user-cloud-share.wake.com.apple.applemediaservices.multiuser</string>
			</dict>
			<key>com.apple.homed.user-cloud-share.wake.com.apple.applemediaservices.multiuser.qa</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.user-cloud-share.wake.com.apple.applemediaservices.multiuser.qa</string>
			</dict>
			<key>com.apple.language.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.language.changed</string>
			</dict>
			<key>com.apple.mobile.keybagd.lock_status</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.lock_status</string>
			</dict>
			<key>com.apple.os-eligibility-domain.input-needed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.os-eligibility-domain.input-needed</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.amsaccountsd.multiuser</key>
		<true/>
		<key>com.apple.amsprivateidentifiers</key>
		<true/>
		<key>com.apple.aps.amsaccountsd</key>
		<true/>
		<key>com.apple.aps.amsprivateidentifiers</key>
		<true/>
		<key>com.apple.private.alloy.amsaccountsync-idswake</key>
		<true/>
		<key>com.apple.xpc.amsaccountsd</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AppleMediaServices.framework/amsaccountsd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
