# diff/tvOS/System/Library/LaunchDaemons/tmpq46su1hi.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.servicesintelligenced</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.servicesintelligenced.cron</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>86400</integer>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>82800</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresNetworkConnectivity</key>
				<true/>
				<key>RequiresUserInactivity</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.StoreServices.SSAccountStore.activeaccountchanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.StoreServices.SSAccountStore.activeaccountchanged</string>
			</dict>
			<key>com.apple.StoreServices.StorefrontChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.StoreServices.StorefrontChanged</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.servicesintelligenced</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/ServicesIntelligence.framework/servicesintelligenced</string>
	</array>
	<key>StandardErrorPath</key>
	<string>/tmp/servicesintelligenced.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/servicesintelligenced.log</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
