# diff/tvOS/System/Library/LaunchDaemons/tmpszirg6qu.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.userprofilesd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.accounts.changes</key>
		<dict>
			<key>AppleAccountAddedOrDeleted</key>
			<dict>
				<key>com.apple.account.AppleAccount</key>
				<dict>
					<key>ChangeTypes</key>
					<array>
						<integer>1</integer>
						<integer>3</integer>
					</array>
				</dict>
			</dict>
			<key>GameCenterAccountAddedOrDeleted</key>
			<dict>
				<key>com.apple.account.GameCenter</key>
				<dict>
					<key>ChangeTypes</key>
					<array>
						<integer>1</integer>
						<integer>3</integer>
					</array>
				</dict>
			</dict>
			<key>iTunesStoreAccountAddedOrDeleted</key>
			<dict>
				<key>com.apple.account.iTunesStore</key>
				<dict>
					<key>ChangeTypes</key>
					<array>
						<integer>1</integer>
						<integer>3</integer>
					</array>
				</dict>
			</dict>
		</dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.userprofilesd.maintenanceTask</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>69120</integer>
				</dict>
				<key>RequiresNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.authkit.simple-profile-list-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.authkit.simple-profile-list-changed</string>
			</dict>
			<key>com.apple.siri.rmv.shared_user_update</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.rmv.shared_user_update</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.userprofiles</key>
		<true/>
		<key>com.apple.userprofiles</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/UserProfilesDaemon.framework/Support/userprofilesd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
