# diff/tvOS/System/Library/LaunchDaemons/tmps8bd25q_.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.SiriRequestService</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>CKPersistentStorageKeysAndValuesUpdated.SiriUserCorrections</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysAndValuesUpdated.SiriUserCorrections</string>
			</dict>
			<key>CKPersistentStorageKeysAndValuesUpdated.synched</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysAndValuesUpdated.synched</string>
			</dict>
			<key>CKPersistentStorageKeysRemoved.SiriUserCorrections</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysRemoved.SiriUserCorrections</string>
			</dict>
			<key>CKPersistentStorageKeysRemoved.synched</key>
			<dict>
				<key>Name</key>
				<string>CKPersistentStorageKeysRemoved.synched</string>
			</dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>AppleDatePreferencesChangedNotification</key>
			<dict>
				<key>Notification</key>
				<string>AppleDatePreferencesChangedNotification</string>
			</dict>
			<key>com.apple.MobileAsset.VoiceServices.VoiceResources.new-asset-installed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MobileAsset.VoiceServices.VoiceResources.new-asset-installed</string>
			</dict>
			<key>com.apple.VideoSubscriberAccount.DidRegisterSubscription</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.VideoSubscriberAccount.DidRegisterSubscription</string>
			</dict>
			<key>com.apple.WatchListKit.WLKStoredConfigurationDidChangeDispatchNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.WatchListKit.WLKStoredConfigurationDidChangeDispatchNotification</string>
			</dict>
			<key>com.apple.homed.user-cloud-share.wake.com.apple.siri.data</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.user-cloud-share.wake.com.apple.siri.data</string>
			</dict>
			<key>com.apple.homed.user-cloud-share.wake.com.apple.siri.zonesharing</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.homed.user-cloud-share.wake.com.apple.siri.zonesharing</string>
			</dict>
			<key>com.apple.managedconfiguration.effectivesettingschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.effectivesettingschanged</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
			<key>com.apple.siri.preheat.quiet</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.preheat.quiet</string>
			</dict>
			<key>com.apple.voiceservices.notification.voice-update</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voiceservices.notification.voice-update</string>
			</dict>
			<key>com.apple.voicetrigger.EarlyDetect</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.voicetrigger.EarlyDetect</string>
			</dict>
		</dict>
		<key>com.apple.rapport.matching</key>
		<dict>
			<key>WakeUpRequest</key>
			<dict>
				<key>serviceType</key>
				<string>com.apple.siri.wakeup</string>
				<key>type</key>
				<string>server</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.assistant.analytics-observation.remora</key>
		<true/>
		<key>com.apple.assistant.audio-session-assertion.remora</key>
		<true/>
		<key>com.apple.assistant.client.remora</key>
		<true/>
		<key>com.apple.assistant.device-context.remora</key>
		<true/>
		<key>com.apple.assistant.dictation.remora</key>
		<true/>
		<key>com.apple.assistant.flow_service.remora</key>
		<true/>
		<key>com.apple.assistant.multiuser.service.remora</key>
		<true/>
		<key>com.apple.assistant.multiuser.update.service.remora</key>
		<true/>
		<key>com.apple.assistant.security.remora</key>
		<true/>
		<key>com.apple.assistant.settings.remora</key>
		<true/>
		<key>com.apple.assistant.sync.remora</key>
		<true/>
		<key>com.apple.assistantd.managedstorage.remora</key>
		<true/>
		<key>com.apple.private.alloy.siri.device-idswake.remora</key>
		<true/>
		<key>com.apple.private.alloy.siri.icloud-idswake.remora</key>
		<true/>
		<key>com.apple.private.alloy.siri.location-idswake.remora</key>
		<true/>
		<key>com.apple.private.alloy.siri.proxy-idswake.remora</key>
		<true/>
		<key>com.apple.siri.client_lite.remora</key>
		<true/>
		<key>com.apple.siri.device_resolution.remora</key>
		<true/>
		<key>com.apple.siri.external_request.remora</key>
		<true/>
		<key>com.apple.siri.pairedsync.remora</key>
		<true/>
		<key>com.apple.siri.tether.remora</key>
		<true/>
		<key>com.apple.siri.vocabularyupdates.remora</key>
		<true/>
	</dict>
	<key>MultipleInstances</key>
	<string>Managed</string>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/AssistantServices.framework/assistantd</string>
		<string>--remora</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
