# diff/tvOS/System/Library/LaunchDaemons/tmptpo15rxa.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.modelcatalogd</string>
	<key>LaunchEvents</key>
	<dict/>
	<key>MachServices</key>
	<dict>
		<key>com.apple.modelcatalog.catalog</key>
		<true/>
		<key>com.apple.modelcatalog.subscriptions</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/ModelCatalogRuntime.framework/modelcatalogd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
