# diff/tvOS/System/Library/LaunchDaemons/tmptkkc1orm.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.corecaptured</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.corecaptured</key>
			<dict>
				<key>IOProviderClass</key>
				<string>CCCapture</string>
				<key>IORegistryName</key>
				<string>CCCapture</string>
			</dict>
			<key>com.apple.corecaptured-driverkit</key>
			<dict>
				<key>IOClass</key>
				<string>CCCapture</string>
				<key>IOPropertyMatch</key>
				<dict>
					<key>CCIOServiceObjType</key>
					<integer>5001</integer>
				</dict>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.MCX._managementStatusChangedForDomains</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MCX._managementStatusChangedForDomains</string>
			</dict>
			<key>com.apple.ManagedConfiguration.profileListChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.ManagedConfiguration.profileListChanged</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.corecaptured</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/corecaptured</string>
	</array>
	<key>RemoteServices</key>
	<dict>
		<key>com.apple.corecaptured.remoteservice</key>
		<dict>
			<key>RequireEntitlement</key>
			<string>com.apple.corecaptured.remoteservice-access</string>
			<key>ServiceProperties</key>
			<dict>
				<key>UsesRemoteXPC</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>StartInterval</key>
	<integer>86400</integer>
</dict>
</plist>
