# diff/tvOS/System/Library/LaunchDaemons/tmpublg6x8a.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.photoanalysisd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.photoanalysisd.backgroundanalysis</key>
			<dict>
				<key>Dependencies</key>
				<array>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.mediaanalysisd.photos.face</string>
						<key>DependencyType</key>
						<string>DependencyTypeCompletion</string>
					</dict>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.mediaanalysisd.photos.face.fullResults</string>
						<key>DependencyType</key>
						<string>DependencyTypeResult</string>
						<key>ResultDependencyBatchSize</key>
						<string>500</string>
					</dict>
				</array>
				<key>ExpectedDuration</key>
				<integer>2700</integer>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.photoanalysisd.sequentialProcessing</string>
				<key>PowerNap</key>
				<true/>
				<key>PreventsDeviceSleep</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>ProducedResultIdentifiers</key>
				<array>
					<string>com.apple.photoanalysisd.backgroundanalysis</string>
				</array>
				<key>RandomInitialDelay</key>
				<integer>2700</integer>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>21600</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>6000</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
			</dict>
			<key>com.apple.photoanalysisd.featuredcontent</key>
			<dict>
				<key>Dependencies</key>
				<array>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.photoanalysisd.backgroundanalysis</string>
						<key>DependencyType</key>
						<string>DependencyTypeCompletion</string>
					</dict>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.photoanalysisd.highlightenrichment</string>
						<key>DependencyType</key>
						<string>DependencyTypeCompletion</string>
					</dict>
				</array>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.photoanalysisd.sequentialProcessing</string>
				<key>PowerNap</key>
				<true/>
				<key>PreventsDeviceSleep</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>ProducedResultIdentifiers</key>
				<array>
					<string>com.apple.photoanalysisd.featuredcontent</string>
				</array>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>6000</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
			</dict>
			<key>com.apple.photoanalysisd.graphNonIntensiveTasks</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.photoanalysisd.sequentialProcessing</string>
				<key>PowerNap</key>
				<true/>
				<key>PreventsDeviceSleep</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>7200</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>6000</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.photoanalysisd.highlightenrichment</key>
			<dict>
				<key>Dependencies</key>
				<array>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.photoanalysisd.backgroundanalysis</string>
						<key>DependencyType</key>
						<string>DependencyTypeCompletion</string>
					</dict>
				</array>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.photoanalysisd.sequentialProcessing</string>
				<key>PowerNap</key>
				<true/>
				<key>PreventsDeviceSleep</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>ProducedResultIdentifiers</key>
				<array>
					<string>com.apple.photoanalysisd.highlightenrichment</string>
				</array>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>21600</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>6000</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
			</dict>
			<key>com.apple.photoanalysisd.internal</key>
			<dict>
				<key>Dependencies</key>
				<array>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.photoanalysisd.backgroundanalysis</string>
						<key>DependencyType</key>
						<string>DependencyTypeCompletion</string>
					</dict>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.photoanalysisd.highlightenrichment</string>
						<key>DependencyType</key>
						<string>DependencyTypeCompletion</string>
					</dict>
				</array>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.photoanalysisd.sequentialProcessing</string>
				<key>PowerNap</key>
				<true/>
				<key>PreventsDeviceSleep</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
			</dict>
			<key>com.apple.photoanalysisd.libraryprocessing</key>
			<dict>
				<key>Dependencies</key>
				<array>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.photoanalysisd.backgroundanalysis</string>
						<key>DependencyType</key>
						<string>DependencyTypeCompletion</string>
					</dict>
				</array>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.photoanalysisd.sequentialProcessing</string>
				<key>PowerNap</key>
				<true/>
				<key>PreventsDeviceSleep</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>6000</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
			</dict>
			<key>com.apple.photoanalysisd.music</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.photoanalysisd.sequentialProcessing</string>
				<key>PowerNap</key>
				<true/>
				<key>PreventsDeviceSleep</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RandomInitialDelay</key>
				<integer>7200</integer>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>28800</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>23040</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.photoanalysisd.photoLibraryOnBatteryTasks</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.photoanalysisd.sequentialProcessing</string>
				<key>PowerNap</key>
				<true/>
				<key>PreventsDeviceSleep</key>
				<false/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>28800</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>23040</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.photoanalysisd.photolibraryNonIntensiveTasks</key>
			<dict>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.photoanalysisd.sequentialProcessing</string>
				<key>PowerNap</key>
				<true/>
				<key>PreventsDeviceSleep</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>7200</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>6000</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
			<key>com.apple.photoanalysisd.searchdonation</key>
			<dict>
				<key>Dependencies</key>
				<array>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.photoanalysisd.backgroundanalysis</string>
						<key>DependencyType</key>
						<string>DependencyTypeCompletion</string>
					</dict>
				</array>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.photoanalysisd.sequentialProcessing</string>
				<key>PowerNap</key>
				<true/>
				<key>PreventsDeviceSleep</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>21600</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>6000</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
			</dict>
			<key>com.apple.photoanalysisd.wallpaper</key>
			<dict>
				<key>Dependencies</key>
				<array>
					<dict>
						<key>DependencyIdentifier</key>
						<string>com.apple.photoanalysisd.backgroundanalysis</string>
						<key>DependencyType</key>
						<string>DependencyTypeCompletion</string>
					</dict>
				</array>
				<key>GroupConcurrencyLimit</key>
				<integer>1</integer>
				<key>GroupName</key>
				<string>com.apple.photoanalysisd.sequentialProcessing</string>
				<key>PowerNap</key>
				<true/>
				<key>PreventsDeviceSleep</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
					<key>MinDurationBetweenInstances</key>
					<integer>6000</integer>
				</dict>
				<key>RequiresBuddyComplete</key>
				<true/>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>RTLocationsOfInterestDidClearNotification</key>
			<dict>
				<key>Notification</key>
				<string>RTLocationsOfInterestDidClearNotification</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.photoanalysisd</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/PhotoAnalysis.framework/Support/photoanalysisd</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
