# diff/tvOS/System/Library/LaunchDaemons/tmpvp9sgfq0.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.passd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.appprotection.change</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appprotection.change</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
		</dict>
		<key>com.apple.nsurlsessiond-events</key>
		<dict>
			<key>LaunchOnDemand</key>
			<dict>
				<key>ClientIdentifier</key>
				<string>com.apple.passd</string>
			</dict>
		</dict>
		<key>com.apple.parsecd.pushevents</key>
		<dict>
			<key>flightPushSubscription</key>
			<dict>
				<key>SubscriptionName</key>
				<string>flightsPush</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.passd.account</key>
		<true/>
		<key>com.apple.passd.aps</key>
		<true/>
		<key>com.apple.passd.assertions</key>
		<true/>
		<key>com.apple.passd.bulletins</key>
		<true/>
		<key>com.apple.passd.cloud-store</key>
		<true/>
		<key>com.apple.passd.device-registration</key>
		<true/>
		<key>com.apple.passd.discovery</key>
		<true/>
		<key>com.apple.passd.in-app-payment</key>
		<true/>
		<key>com.apple.passd.library</key>
		<true/>
		<key>com.apple.passd.nf-events</key>
		<true/>
		<key>com.apple.passd.payment</key>
		<true/>
		<key>com.apple.passd.payment-continuity</key>
		<true/>
		<key>com.apple.passd.peer-payment</key>
		<true/>
		<key>com.apple.passd.search</key>
		<true/>
		<key>com.apple.passd.sharing-channel</key>
		<true/>
		<key>com.apple.passd.trusted-device-enrollment-info-provider</key>
		<true/>
		<key>com.apple.passd.usage</key>
		<dict>
			<key>_ExperimentalNonLaunching</key>
			<true/>
		</dict>
		<key>com.apple.private.alloy.applepay-idswake</key>
		<true/>
		<key>com.apple.private.alloy.applepay.sharing-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.NanoPassbook</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.Passbook</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.PassbookUIService</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/PassKitCore.framework/passd</string>
	</array>
	<key>PublishesEvents</key>
	<array>
		<string>com.apple.passd.matching</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
