# diff/tvOS/System/Library/LaunchDaemons/tmpxumcq02x.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>OS_ACTIVITY_ENABLE_DYNAMIC</key>
		<string>1</string>
	</dict>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamMemoryLimit</key>
		<integer>7500</integer>
		<key>JetsamPriority</key>
		<integer>-49</integer>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.apple.locationd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>CLPrivacyManager.B</key>
			<dict>
				<key>Interval</key>
				<integer>604800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
			<key>CLPrivacyManager.C</key>
			<dict>
				<key>Interval</key>
				<integer>604800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
			<key>CLPrivacyManager.D</key>
			<dict>
				<key>Interval</key>
				<integer>604800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RequiresClassC</key>
				<false/>
			</dict>
			<key>com.apple.locationd.CLClientManager.coreanalytics</key>
			<dict>
				<key>GracePeriod</key>
				<integer>86400</integer>
				<key>Interval</key>
				<integer>604800</integer>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.eed</key>
		<true/>
		<key>com.apple.aps.locationd.coldstorage.cloudkit</key>
		<true/>
		<key>com.apple.aps.locationd.vo2max.cloudkit</key>
		<true/>
		<key>com.apple.locationauth.pairedsync</key>
		<true/>
		<key>com.apple.locationd.diagnostic</key>
		<true/>
		<key>com.apple.locationd.eedhelper</key>
		<true/>
		<key>com.apple.locationd.gathering</key>
		<true/>
		<key>com.apple.locationd.harvester</key>
		<true/>
		<key>com.apple.locationd.pairedsync</key>
		<true/>
		<key>com.apple.locationd.registration</key>
		<true/>
		<key>com.apple.locationd.routine</key>
		<true/>
		<key>com.apple.locationd.simulation</key>
		<true/>
		<key>com.apple.locationd.smoother</key>
		<true/>
		<key>com.apple.locationd.spi</key>
		<true/>
		<key>com.apple.locationd.synchronous</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.findmy</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.locationd.LocationNotificationBundle</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/locationd</string>
	</array>
	<key>PublishesEvents</key>
	<dict>
		<key>com.apple.locationd.nomovement.stateupdate</key>
		<dict>
			<key>RequireTrusted</key>
			<true/>
		</dict>
		<key>com.apple.locationd.workout.workoutStateStopDetectionOrAutoPauseUpdate</key>
		<dict>
			<key>RequireTrusted</key>
			<true/>
		</dict>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>86400</integer>
	<key>ThrottleInterval</key>
	<integer>3</integer>
</dict>
</plist>
