# diff/tvOS/System/Library/LaunchDaemons/tmpzr65i9w2.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.apple.atc</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.Music-AllowsCellularDataDownloads</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.Music-AllowsCellularDataDownloads</string>
			</dict>
			<key>com.apple.Music-keeplocalstatuschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.atc.xpc.runkeeplocaltask</string>
			</dict>
			<key>com.apple.datamigrator.datamigrationcompletecontinuerestore</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.datamigrator.datamigrationcompletecontinuerestore</string>
			</dict>
			<key>com.apple.fairplayd.resync-fpkeybag</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.fairplayd.resync-fpkeybag</string>
			</dict>
			<key>com.apple.itunescloudd.artworkDownloadsDidCompleteNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.itunescloudd.artworkDownloadsDidCompleteNotification</string>
			</dict>
			<key>com.apple.itunesstored.accountschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.itunesstored.accountschanged</string>
			</dict>
			<key>com.apple.mobileipod-prefsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobileipod-prefsChanged</string>
			</dict>
			<key>com.apple.mobileipod.displayvalueschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobileipod.displayvalueschanged</string>
			</dict>
			<key>com.apple.mobileipod.keeplocalstatechanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobileipod.keeplocalstatechanged</string>
			</dict>
			<key>com.apple.mobileipod.librarychanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobileipod.librarychanged</string>
			</dict>
			<key>com.apple.mobileipod.noncontentspropertieschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobileipod.noncontentspropertieschanged</string>
			</dict>
			<key>com.apple.nanomusic.sync.defaults</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanomusic.sync.defaults</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidpair</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidunpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidunpair</string>
			</dict>
			<key>com.apple.nanoregistry.paireddevicedidchangeversion</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.paireddevicedidchangeversion</string>
			</dict>
			<key>com.apple.springboard.pluggedin</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.pluggedin</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.atc</key>
		<true/>
		<key>com.apple.atc.xpc</key>
		<true/>
		<key>com.apple.atc.xpc.downloadprogress</key>
		<true/>
		<key>com.apple.atc.xpc.icappobserverdelegate</key>
		<true/>
		<key>com.apple.atc.xpc.observer</key>
		<true/>
		<key>com.apple.atc.xpc.runkeeplocaltask</key>
		<true/>
		<key>com.apple.atc.xpc.service</key>
		<true/>
		<key>com.apple.atc.xpc.sessions</key>
		<true/>
		<key>com.apple.atc2</key>
		<true/>
		<key>com.apple.private.alloy.airtraffic-idswake</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/atc</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>mobile</string>
	<key>com.apple.security.exception.mach-lookup.global-name</key>
	<array>
		<string>com.apple.itunescloudd.playactivity</string>
	</array>
</dict>
</plist>
