# beta/iPadOS/System/Library/LaunchDaemons/com.apple.findmy.findmylocated.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.findmy.findmylocated</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>NILocalDeviceStartedInteractingWithTokenNotification</key>
			<dict>
				<key>Notification</key>
				<string>NILocalDeviceStartedInteractingWithTokenNotification</string>
			</dict>
			<key>SPAccountRemovedNotification</key>
			<dict>
				<key>Notification</key>
				<string>SPAccountRemovedNotification</string>
			</dict>
			<key>com.apple.locationd.authorization</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.authorization</string>
			</dict>
			<key>com.apple.mobile.keybagd.first_unlock</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.first_unlock</string>
			</dict>
			<key>com.apple.nanoregistry.devicedidunpair</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.devicedidunpair</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.findmy.findmylocated.post-install</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>PostInstall</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<false/>
				<key>RequireNetworkConnectivity</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.findmy.findmylocate.fenceservice</key>
		<true/>
		<key>com.apple.findmy.findmylocate.friendshipservice</key>
		<true/>
		<key>com.apple.findmy.findmylocate.locationservice</key>
		<true/>
		<key>com.apple.findmy.findmylocate.settings</key>
		<true/>
		<key>com.apple.findmy.push.com.apple.findmy.container</key>
		<true/>
		<key>com.apple.icloud.fmfd.aps</key>
		<true/>
		<key>com.apple.private.alloy.fmd-idswake</key>
		<true/>
		<key>com.apple.private.alloy.fmf-idswake</key>
		<true/>
		<key>com.apple.private.alloy.fmf.local-idswake</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/findmylocated</string>
	<key>Umask</key>
	<integer>23</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
