# beta/iPadOS/System/Library/LaunchDaemons/com.apple.textcomposerd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.textcomposerd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.generativeexperiences.externaltextcomposition</key>
		<dict>
			<key>#Else</key>
			<false/>
			<key>#IfFeatureFlagEnabled</key>
			<string>TextComposer/TextComposerDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.generativeexperiences.textcomposition</key>
		<dict>
			<key>#Else</key>
			<false/>
			<key>#IfFeatureFlagEnabled</key>
			<string>TextComposer/TextComposerDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
	</dict>
	<key>NSAppTransportSecurity</key>
	<dict>
		<key>NSAllowsArbitraryLoads</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/textcomposerd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
