# beta/iPadOS/System/Library/LaunchDaemons/com.apple.SpringBoard.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>CFBundleIdentifier</key>
	<string>com.apple.springboard</string>
	<key>EmbeddedHomeScreen</key>
	<true/>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CA_ASSERT_MAIN_THREAD_TRANSACTIONS</key>
		<string>1</string>
	</dict>
	<key>HighPriorityIO</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.SpringBoard</string>
	<key>LimitLoadFromHardware</key>
	<dict>
		<key>osenvironment</key>
		<array>
			<string>device-recovery</string>
			<string>diagnostics</string>
		</array>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>PurpleSystemAppPort</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.DragUI.druid.system</key>
		<true/>
		<key>com.apple.SBUserNotification</key>
		<true/>
		<key>com.apple.UIKit.KeyboardManagement.hosted</key>
		<true/>
		<key>com.apple.UIKit.statusbarserver</key>
		<true/>
		<key>com.apple.UserNotifications.OneTimeCodeService</key>
		<true/>
		<key>com.apple.aps.usernotifications</key>
		<true/>
		<key>com.apple.aps.usernotifications.dev</key>
		<true/>
		<key>com.apple.bulletinboard.dataproviderconnection</key>
		<true/>
		<key>com.apple.bulletinboard.observerconnection</key>
		<true/>
		<key>com.apple.bulletinboard.serverconduitconnection</key>
		<true/>
		<key>com.apple.bulletinboard.settingsconnection</key>
		<true/>
		<key>com.apple.controlcenter.remoteservice</key>
		<true/>
		<key>com.apple.corespotlight.daemon.notifications</key>
		<true/>
		<key>com.apple.frontboard.systemappservices</key>
		<true/>
		<key>com.apple.hangtracermonitor</key>
		<true/>
		<key>com.apple.notificationcenter.widgetcontrollerconnection</key>
		<true/>
		<key>com.apple.pegasus</key>
		<true/>
		<key>com.apple.pencilkit.pktouchdetection</key>
		<true/>
		<key>com.apple.posterboardui.services</key>
		<true/>
		<key>com.apple.proactive.HomeScreenSuggestion.xpc</key>
		<true/>
		<key>com.apple.proactive.UserEducationSuggestion.client-listener.xpc</key>
		<true/>
		<key>com.apple.siri.activation</key>
		<true/>
		<key>com.apple.siri.activation.service</key>
		<true/>
		<key>com.apple.siri.assessment-mode-restriction</key>
		<true/>
		<key>com.apple.springboard.backgroundappservices</key>
		<dict>
			<key>HideUntilCheckIn</key>
			<true/>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.springboard.blockableservices</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.springboard.carditemscontroller</key>
		<true/>
		<key>com.apple.springboard.oswatchdog</key>
		<true/>
		<key>com.apple.springboard.services</key>
		<dict>
			<key>HideUntilCheckIn</key>
			<true/>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.springboard.statusbarservices</key>
		<dict/>
		<key>com.apple.systemstatus</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.systemstatus.activityattribution</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.systemstatus.publisher</key>
		<dict>
			<key>#IfFeatureFlagDisabled</key>
			<string>SystemStatus/EmbeddedDaemon</string>
			<key>#Then</key>
			<true/>
		</dict>
		<key>com.apple.usernotifications.delegate.com.apple.donotdisturb</key>
		<true/>
		<key>com.apple.usernotifications.remotenotificationservice</key>
		<true/>
		<key>com.apple.usernotifications.systemservice</key>
		<true/>
		<key>com.apple.usernotifications.usernotificationservice</key>
		<true/>
		<key>com.apple.usernotifications.usernotificationsettingsservice</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/SpringBoard.app/SpringBoard</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>UserName</key>
	<string>mobile</string>
	<key>_AdditionalProperties</key>
	<dict>
		<key>RunningBoard</key>
		<dict>
			<key>UnderlyingAssertion</key>
			<string>UnderlyingDarwinRoleUI</string>
		</dict>
	</dict>
	<key>_Conclave</key>
	<string>com.apple.springboard.conclave</string>
	<key>_ForceHardMemoryLimit</key>
	<dict>
		<key>#IfFeatureFlagEnabled</key>
		<string>MemoryTelemetry/SelectProcessesHardMemoryLimitsInternal</string>
		<key>#Then</key>
		<true/>
	</dict>
	<key>_LimitLoadFromClarityMode</key>
	<true/>
</dict>
</plist>
