# beta/iPadOS/System/Library/LaunchDaemons/com.apple.device-o-matic.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.device-o-matic</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.device-o-matic.billboard</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOProviderClass</key>
				<string>AppleUSBHostBillboardDevice</string>
			</dict>
			<key>com.apple.device-o-matic.overcurrent</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyExistsMatch</key>
				<string>OvercurrentCount</string>
			</dict>
			<key>com.apple.device-o-matic.usbiokitlaunch</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyMatch</key>
				<dict>
					<key>removable</key>
					<true/>
				</dict>
				<key>IOProviderClass</key>
				<string>IOUSBHostDevice</string>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/System/Library/CoreServices/DeviceOMatic.app/DeviceOMatic</string>
</dict>
</plist>
