# beta/iPadOS/System/Library/LaunchDaemons/com.apple.mediaartworkd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.mediaartworkd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mediaartworkd.xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>AdaptiveUtility</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/MediaServices.framework/Support/mediaartworkd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
