# beta/iPadOS/System/Library/LaunchDaemons/com.apple.liveactivitiesd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>MallocStackLoggingDirectory</key>
		<string>/private/var/tmp/com.apple.sessions</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>Crashed</key>
		<true/>
	</dict>
	<key>Label</key>
	<string>com.apple.liveactivitiesd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.alarm</key>
		<true/>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>application-installed</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>application-uninstalled</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>accessory-attached</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.accpowersources.attach</string>
			</dict>
			<key>accessory-power-source-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.accpowersources.source</string>
			</dict>
			<key>com.apple.appprotection.change</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.appprotection.change</string>
			</dict>
			<key>com.apple.liveactivitiesd.replicatorParticipant.message</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.liveactivitiesd.replicatorParticipant.message</string>
			</dict>
			<key>com.apple.liveactivitiesd.replicatorParticipant.record</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.liveactivitiesd.replicatorParticipant.record</string>
			</dict>
			<key>keybag-status-changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.lock_status</string>
			</dict>
			<key>ui-unlocked</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.lockstate</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.sessioncore.PushNotifications</key>
		<true/>
		<key>com.apple.aps.sessioncore.PushNotifications.dev</key>
		<true/>
		<key>com.apple.private.alloy.sessionkit-idswake</key>
		<true/>
		<key>com.apple.sessionservices</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/SessionCore.framework/Support/liveactivitiesd</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>ThrottleInterval</key>
	<string>3</string>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
