# beta/iPadOS/System/Library/LaunchDaemons/com.apple.intelligenceplatformd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AppIntents</key>
	<dict>
		<key>MetadataAbsolutePaths</key>
		<array>
			<string>/System/Library/PrivateFrameworks/IntelligencePlatformCore.framework</string>
		</array>
	</dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>Label</key>
	<string>com.apple.intelligenceplatformd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.biome.compute.publisher</key>
		<dict>
			<key>com.apple.intelligenceplatformd.EntityTaggingSubscription</key>
			<dict>
				<key>client</key>
				<string>com.apple.intelligenceplatformd</string>
				<key>graph</key>
				<dict>
					<key>keyPath</key>
					<string>eventBody.viewName</string>
					<key>name</key>
					<string>filter</string>
					<key>upstream</key>
					<dict>
						<key>name</key>
						<string>streamPublisher</string>
						<key>streamIdentifier</key>
						<string>IntelligencePlatform.Views.Updated</string>
						<key>version</key>
						<integer>1</integer>
					</dict>
					<key>value</key>
					<string>entityTagging</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.intelligenceplatformd.EntityTaggingSubscription</string>
			</dict>
			<key>com.apple.intelligenceplatformd.MotionSegment</key>
			<dict>
				<key>client</key>
				<string>com.apple.intelligenceplatformd</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>Motion.TripSegment</string>
					<key>streamType</key>
					<string>2</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.intelligenceplatformd.MotionSegment</string>
			</dict>
			<key>com.apple.intelligenceplatformd.NowPlaying</key>
			<dict>
				<key>client</key>
				<string>com.apple.intelligenceplatformd</string>
				<key>graph</key>
				<dict>
					<key>keyPath</key>
					<string>eventBody.bundleID</string>
					<key>name</key>
					<string>filter</string>
					<key>upstream</key>
					<dict>
						<key>name</key>
						<string>streamPublisher</string>
						<key>streamIdentifier</key>
						<string>NowPlaying</string>
						<key>version</key>
						<integer>1</integer>
					</dict>
					<key>value</key>
					<string>com.apple.Music</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.intelligenceplatformd.NowPlaying</string>
			</dict>
			<key>com.apple.intelligenceplatformd.walletClassicOrder</key>
			<dict>
				<key>client</key>
				<string>com.apple.intelligenceplatformd</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>WalletPaymentsCommerce.FoundIn.ClassicOrder</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.intelligenceplatformd.walletClassicOrder</string>
			</dict>
			<key>com.apple.intelligenceplatformd.walletEmailOrder</key>
			<dict>
				<key>client</key>
				<string>com.apple.intelligenceplatformd</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>WalletPaymentsCommerce.FoundIn.OrderEmail</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.intelligenceplatformd.walletEmailOrder</string>
			</dict>
			<key>com.apple.intelligenceplatformd.walletTransactionOrder</key>
			<dict>
				<key>client</key>
				<string>com.apple.intelligenceplatformd</string>
				<key>graph</key>
				<dict>
					<key>name</key>
					<string>streamPublisher</string>
					<key>streamIdentifier</key>
					<string>WalletPaymentsCommerce.FoundIn.Transaction</string>
					<key>version</key>
					<integer>1</integer>
				</dict>
				<key>identifier</key>
				<string>com.apple.intelligenceplatformd.walletTransactionOrder</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>ABAddressBookMeCardChangeDistributedNotification</key>
			<dict>
				<key>Name</key>
				<string>ABAddressBookMeCardChangeDistributedNotification</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>RTLocationsOfInterestDidChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>RTLocationsOfInterestDidChangeNotification</string>
			</dict>
			<key>__ABDataBaseChangedByOtherProcessNotification</key>
			<dict>
				<key>Notification</key>
				<string>__ABDataBaseChangedByOtherProcessNotification</string>
			</dict>
			<key>com.apple.duetexpertd.sportsTeamsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.duetexpertd.sportsTeamsChanged</string>
			</dict>
			<key>com.apple.intelligenceplatform.StorageSystem.Recovered</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.intelligenceplatform.StorageSystem.Recovered</string>
			</dict>
			<key>com.apple.photos.DidUpdateAutonamingUserFeedback</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.photos.DidUpdateAutonamingUserFeedback</string>
			</dict>
			<key>com.apple.siri.uaf.com.apple.siri.understanding</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.siri.uaf.com.apple.siri.understanding</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.AssetRegistryOverrideStoreWeeklyCleanup</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>CPUIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>604800</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.AutonamingMetricsLogging</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>CPUIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.EventAndFeedbackDailyProcessing</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>CPUIntensive</key>
				<false/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.EventStreamsDataCollection</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>CPUIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.GlobalKnowledgeTwoHourPrune</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>CPUIntensive</key>
				<false/>
				<key>Interval</key>
				<integer>7200</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
			<key>com.apple.intelligenceplatform.IntelligencePlatformCore.LongitudinalDataCollection</key>
			<dict>
				<key>ActivityGroupName</key>
				<string>com.apple.intelligenceplatform.IntelligencePlatformCore</string>
				<key>CPUIntensive</key>
				<true/>
				<key>Interval</key>
				<integer>86400</integer>
				<key>MaxConcurrencyLimit</key>
				<integer>1</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.intelligenceplatform.AssetRegistry</key>
		<true/>
		<key>com.apple.intelligenceplatform.EntityResolution</key>
		<true/>
		<key>com.apple.intelligenceplatform.EventLog</key>
		<true/>
		<key>com.apple.intelligenceplatform.Feedback</key>
		<true/>
		<key>com.apple.intelligenceplatform.InferenceSupport</key>
		<true/>
		<key>com.apple.intelligenceplatform.Internal</key>
		<true/>
		<key>com.apple.intelligenceplatform.KGSimulation</key>
		<true/>
		<key>com.apple.intelligenceplatform.Knosis</key>
		<true/>
		<key>com.apple.intelligenceplatform.Lighthouse</key>
		<true/>
		<key>com.apple.intelligenceplatform.Sysdiagnose</key>
		<true/>
		<key>com.apple.private.appintents.delegate.com.apple.intelligenceplatformd</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/IntelligencePlatformCore.framework/intelligenceplatformd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
