# beta/iPadOS/System/Library/LaunchDaemons/com.apple.storagekitd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.storagekitd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>com.apple.storagekit.resize-match</key>
			<dict>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPathMatch</key>
				<string>IODeviceTree:/defaults</string>
				<key>IOPropertyExistsMatch</key>
				<string>sk-disk-resize</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.storagekitd</key>
		<true/>
		<key>com.apple.storagekitd.dm</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/storagekitd</string>
	</array>
</dict>
</plist>
