# beta/iPadOS/System/Library/LaunchDaemons/com.apple.pipelined.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>ASL_QUOTA_DISABLED</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.pipelined</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.nsurlsessiond-events</key>
		<dict>
			<key>LaunchOnDemand</key>
			<dict>
				<key>ClientIdentifier</key>
				<string>com.apple.pipelined</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.indoor.ClearTiles</key>
			<dict>
				<key>Interval</key>
				<real>43200</real>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
			<key>com.apple.indoor.VacuumDBs</key>
			<dict>
				<key>Interval</key>
				<integer>518400</integer>
				<key>Priority</key>
				<string>Maintenance</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.pipelined</key>
		<true/>
		<key>com.apple.pipelined.maintenance</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/pipelined</string>
	</array>
</dict>
</plist>
