# beta/iPadOS/System/Library/LaunchDaemons/com.apple.iapauthd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.iapauthd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.iapauthd</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.iapauthd.xpc</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/IAPAuthentication.framework/Support/iapauthd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
