# beta/iPadOS/System/Library/LaunchDaemons/com.apple.ioupsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>Label</key>
	<string>com.apple.ioupsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>HIDDevice attached</key>
			<dict>
				<key>DeviceUsagePairs</key>
				<array>
					<dict>
						<key>DeviceUsagePage</key>
						<integer>132</integer>
					</dict>
					<dict>
						<key>DeviceUsagePage</key>
						<integer>133</integer>
					</dict>
					<dict>
						<key>DeviceUsage</key>
						<integer>20</integer>
						<key>DeviceUsagePage</key>
						<integer>65280</integer>
					</dict>
					<dict>
						<key>DeviceUsage</key>
						<integer>6</integer>
						<key>DeviceUsagePage</key>
						<integer>132</integer>
					</dict>
				</array>
				<key>IOMatchLaunchStream</key>
				<false/>
				<key>IOProviderClass</key>
				<string>IOHIDDevice</string>
			</dict>
			<key>SMCAccessory attached</key>
			<dict>
				<key>IOMatchAll</key>
				<true/>
				<key>IOMatchLaunchStream</key>
				<true/>
				<key>IOPropertyExistsMatch</key>
				<string>AccessoryDetails</string>
				<key>IOProviderClass</key>
				<string>IOPMPowerSource</string>
			</dict>
		</dict>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/usr/libexec/ioupsd</string>
</dict>
</plist>
