# beta/iPadOS/System/Library/LaunchDaemons/com.apple.itunescloudd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.itunescloudd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.CTTelephonyCenter</key>
		<dict>
			<key>kCTSettingPhoneNumberChangedNotification</key>
			<dict>
				<key>Name</key>
				<string>kCTSettingPhoneNumberChangedNotification</string>
			</dict>
		</dict>
		<key>com.apple.distnoted.matching</key>
		<dict>
			<key>com.apple.LaunchServices.applicationRegistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationRegistered</string>
			</dict>
			<key>com.apple.LaunchServices.applicationUnregistered</key>
			<dict>
				<key>Name</key>
				<string>com.apple.LaunchServices.applicationUnregistered</string>
			</dict>
			<key>com.apple.StoreServices.SSVFamilyContentDeletionNotification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.StoreServices.SSVFamilyContentDeletionNotification</string>
			</dict>
			<key>com.apple.itunescloud.ICPlayActivityFlushInternalTestRequestNotification</key>
			<dict>
				<key>Name</key>
				<string>com.apple.itunescloud.ICPlayActivityFlushInternalTestRequestNotification</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.datamigrator.datamigrationcompletecontinuerestore</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.datamigrator.datamigrationcompletecontinuerestore</string>
			</dict>
			<key>com.apple.family.family_updated</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.family.family_updated</string>
			</dict>
			<key>com.apple.itunesstored.accountschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.itunesstored.accountschanged</string>
			</dict>
			<key>com.apple.itunesstored.autodownloaddefaultschange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.itunesstored.autodownloaddefaultschange</string>
			</dict>
			<key>com.apple.managedconfiguration.restrictionchanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.managedconfiguration.restrictionchanged</string>
			</dict>
			<key>com.apple.mobileipod.libraryimportdidfinish</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobileipod.libraryimportdidfinish</string>
			</dict>
			<key>com.apple.purplebuddy.setupdone</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.purplebuddy.setupdone</string>
			</dict>
			<key>com.apple.tcc.access.changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.tcc.access.changed</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.itunescloudd</key>
		<true/>
		<key>com.apple.itunescloud.contenttaste</key>
		<true/>
		<key>com.apple.itunescloud.delegate-account-store</key>
		<true/>
		<key>com.apple.itunescloud.delegation-provider</key>
		<true/>
		<key>com.apple.itunescloud.in-app-message-service</key>
		<true/>
		<key>com.apple.itunescloud.library-auth-token-provider</key>
		<true/>
		<key>com.apple.itunescloud.music-subscription-status-service</key>
		<true/>
		<key>com.apple.itunescloud.remote-request-execution-service</key>
		<true/>
		<key>com.apple.itunescloudd.playactivity</key>
		<true/>
		<key>com.apple.itunescloudd.tcchelper</key>
		<true/>
		<key>com.apple.itunescloudd.xpc</key>
		<true/>
		<key>com.apple.private.alloy.itunescloudd-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.Music</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/iTunesCloud.framework/Support/itunescloudd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
