# beta/iPadOS/System/Library/LaunchDaemons/com.apple.DataDetectorsSourceAccess.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>GroupName</key>
	<string>_datadetectors</string>
	<key>Label</key>
	<string>com.apple.DataDetectorsSourceAccess</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.DataDetectorsSourceAccess</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>Program</key>
	<string>/usr/libexec/DataDetectorsSourceAccess</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/DataDetectorsSourceAccess</string>
	</array>
	<key>UserName</key>
	<string>_datadetectors</string>
</dict>
</plist>
