# beta/iPadOS/System/Library/LaunchDaemons/com.apple.lockdownmoded.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.lockdownmoded</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.LockdownMode.accountChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.LockdownMode.accountChanged</string>
			</dict>
			<key>com.apple.kvs.store-did-change.com.apple.LockdownMode</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.kvs.store-did-change.com.apple.LockdownMode</string>
			</dict>
			<key>com.apple.nanoregistry.pairedSync.initialSyncDidComplete</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.pairedSync.initialSyncDidComplete</string>
			</dict>
			<key>com.apple.nanoregistry.paireddevicedidchangeversion</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.paireddevicedidchangeversion</string>
			</dict>
			<key>com.apple.purplebuddy.setupexited</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.purplebuddy.setupexited</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.lockdownmoded.updateManagedConfiguration</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Priority</key>
				<string>Utility</string>
				<key>RandomInitialDelay</key>
				<integer>10</integer>
				<key>Repeating</key>
				<false/>
				<key>RequiresBuddyComplete</key>
				<true/>
			</dict>
			<key>com.apple.lockdownmoded.updateState</key>
			<dict>
				<key>AllowBattery</key>
				<true/>
				<key>Interval</key>
				<integer>3600</integer>
				<key>Priority</key>
				<string>Utility</string>
				<key>Repeating</key>
				<true/>
				<key>RequiresBuddyComplete</key>
				<true/>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.lockdownmoded</key>
		<true/>
		<key>com.apple.private.alloy.lockdownmode-idswake</key>
		<true/>
		<key>com.apple.usernotifications.delegate.com.apple.LockdownMode</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/LockdownMode.framework/lockdownmoded</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
