# beta/iPadOS/System/Library/LaunchDaemons/com.apple.securityresearchdeviceinit.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.apple.securityresearchdeviceinit</string>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/securityresearchdevice-init</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>_LimitLoadToDeviceTree</key>
	<array>
		<string>chosen:research-enabled</string>
		<string>chosen:sec-research-device-erm-enabled</string>
	</array>
</dict>
</plist>
