# beta/iPadOS/System/Library/LaunchDaemons/com.apple.icloudmailagent.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.icloudmailagent</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.accounts.changes</key>
		<dict>
			<key>AppleAccountAddedOrDeleted</key>
			<dict>
				<key>com.apple.account.AppleAccount</key>
				<dict>
					<key>ChangeTypes</key>
					<array>
						<integer>1</integer>
						<integer>3</integer>
					</array>
				</dict>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aps.icloud.mcc</key>
		<true/>
		<key>com.apple.icloudmailagent.secret.xpc</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/icloudmailagent</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
