# beta/iPadOS/System/Library/LaunchDaemons/com.apple.spotlightknowledged.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.spotlightknowledged</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.bg.system.task</key>
		<dict>
			<key>com.apple.corespotlight.knowledge</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<real>86400</real>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
				<key>Resources</key>
				<integer>1</integer>
			</dict>
			<key>com.apple.corespotlight.knowledge.AB</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<real>86400</real>
				</dict>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>A</string>
			</dict>
			<key>com.apple.corespotlight.knowledge.cascadeVerification</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
				<key>Resources</key>
				<string>1</string>
			</dict>
			<key>com.apple.corespotlight.knowledge.indexing</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
				<key>Resources</key>
				<integer>1</integer>
			</dict>
			<key>com.apple.corespotlight.knowledge.inference</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<real>86400</real>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
				<key>Resources</key>
				<integer>1</integer>
			</dict>
			<key>com.apple.corespotlight.knowledge.journals</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
				<key>Resources</key>
				<integer>1</integer>
			</dict>
			<key>com.apple.corespotlight.knowledge.journals.AB</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<real>86400</real>
				</dict>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>A</string>
			</dict>
			<key>com.apple.corespotlight.knowledge.processing</key>
			<dict>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
				<key>ResourceIntensive</key>
				<true/>
				<key>Resources</key>
				<integer>1</integer>
			</dict>
			<key>com.apple.corespotlight.knowledge.report</key>
			<dict>
				<key>Priority</key>
				<string>Utility</string>
				<key>RepeatingTask</key>
				<dict>
					<key>Interval</key>
					<integer>86400</integer>
				</dict>
				<key>RequiresExternalPower</key>
				<true/>
				<key>RequiresNetworkConnectivity</key>
				<false/>
				<key>RequiresProtectionClass</key>
				<string>C</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.CascadeSets.DonateNow</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.CascadeSets.DonateNow</string>
			</dict>
		</dict>
		<key>com.apple.xpc.activity</key>
		<dict>
			<key>com.apple.corespotlight.knowledge.debug.inference</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>2147483647</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<false/>
			</dict>
			<key>com.apple.corespotlight.knowledge.debug.people</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>2147483647</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<false/>
			</dict>
			<key>com.apple.corespotlight.knowledge.debug.reset</key>
			<dict>
				<key>AllowBattery</key>
				<false/>
				<key>Interval</key>
				<integer>2147483647</integer>
				<key>PowerNap</key>
				<true/>
				<key>Priority</key>
				<string>Maintenance</string>
				<key>Repeating</key>
				<false/>
			</dict>
		</dict>
	</dict>
	<key>NO_EnvironmentVariables</key>
	<dict>
		<key>MallocStackLogging</key>
		<integer>1</integer>
		<key>MallocStackLoggingNoCompact</key>
		<integer>1</integer>
		<key>NO_DYLD_INSERT_LIBRARIES</key>
		<string>/usr/lib/libgmalloc.dylib</string>
		<key>NSDebugEnabled</key>
		<integer>1</integer>
		<key>NSZombieEnabled</key>
		<integer>0</integer>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/CoreSpotlight.framework/spotlightknowledged</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>_ExponentialThrottling</key>
	<true/>
</dict>
</plist>
