# beta/iPadOS/System/Library/LaunchDaemons/com.apple.migrationd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>DrainMessagesOnFailedInit</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.migrationd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.migrationkit.internalsettings.dismissed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.migrationkit.internalsettings.dismissed</string>
			</dict>
			<key>com.apple.purplebuddy.setupexited</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.purplebuddy.setupexited</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.migrationd</key>
		<true/>
		<key>com.apple.migrationd.install-coordination</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/MigrationKit.framework/migrationd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
