# diff/iOS/System/Library/NanoLaunchDaemonsAltAccount/com.apple.nanosystemsettingsd.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.nanosystemsettingsd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.coreduetcontext.client_event_stream</key>
		<dict>
			<key>notification</key>
			<dict>
				<key>client-identifier</key>
				<string>com.apple.nanosystemsettingsd.contextstore-registration</string>
			</dict>
		</dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>Active Country Code</key>
			<dict>
				<key>Notification</key>
				<string>CPActiveCountryCodeChanged</string>
			</dict>
			<key>Gizmo Power Update</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.sharing.powerStatusNotification.watch</string>
			</dict>
			<key>Home Country Code</key>
			<dict>
				<key>Notification</key>
				<string>CPHomeCountryCodeChanged</string>
			</dict>
			<key>Keyboard Emojis</key>
			<dict>
				<key>Notification</key>
				<string>AppleKeyboardsEmojiChangedNotification</string>
			</dict>
			<key>Keyboard Language List</key>
			<dict>
				<key>Notification</key>
				<string>AppleKeyboardsPreferencesChangedNotification</string>
			</dict>
			<key>Measurement Unit</key>
			<dict>
				<key>Notification</key>
				<string>AppleMeasurementSystemPreferencesChangedNotification</string>
			</dict>
			<key>NRInitialPairedSyncDidCompleteDarwinNotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanoregistry.pairedSync.initialSyncDidComplete</string>
			</dict>
			<key>Network Country Code</key>
			<dict>
				<key>Notification</key>
				<string>CPNetworkCountryCodeChanged</string>
			</dict>
			<key>Selected Keyboard Language</key>
			<dict>
				<key>Notification</key>
				<string>AppleKeyboardsSettingsChangedNotification</string>
			</dict>
			<key>Temperature Unit</key>
			<dict>
				<key>Notification</key>
				<string>AppleTemperatureUnitPreferencesChangedNotification</string>
			</dict>
			<key>com.apple.gms.availability.notification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.gms.availability.notification</string>
			</dict>
			<key>com.apple.os-eligibility-domain.change.greymatter</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.os-eligibility-domain.change.greymatter</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.nanosystemsettings</key>
		<true/>
		<key>com.apple.nanosystemsettings.pairedsync</key>
		<true/>
		<key>com.apple.private.alloy.systemsettings-idswake</key>
		<true/>
		<key>com.apple.sleep.server</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/NanoSystemSettings.framework/nanosystemsettingsd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
