# beta/iPadOS/System/Library/LaunchAgents/com.apple.ckdiscretionaryd.sync.xml.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CK_IN_SYNC_BUBBLE</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.apple.ckdiscretionaryd.sync</string>
	<key>LimitLoadToSessionType</key>
	<string>PostLogOutSync</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.ckdiscretionaryd</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CloudKitDaemon.framework/Support/ckdiscretionaryd</string>
		<string>--postlogoutsync</string>
	</array>
</dict>
</plist>
